# files/touchegg-2.0.13-1/files
%FILES%
etc/
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/touchegg.desktop
usr/
usr/bin/
usr/bin/touchegg
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/touchegg.service
usr/share/
usr/share/touchegg/
usr/share/touchegg/touchegg.conf
