# files/xen-4.15.0-4/files
%FILES%
boot/
boot/xen-4.15.0arch.config
boot/xen.efi
boot/xen.gz
etc/
etc/bash_completion.d/
etc/bash_completion.d/xl
etc/conf.d/
etc/conf.d/xencommons
etc/conf.d/xendomains
etc/xen/
etc/xen/README
etc/xen/auto/
etc/xen/cpupool
etc/xen/efi-xen.cfg
etc/xen/scripts/
etc/xen/scripts/block
etc/xen/scripts/block-common.sh
etc/xen/scripts/block-drbd-probe
etc/xen/scripts/block-dummy
etc/xen/scripts/block-enbd
etc/xen/scripts/block-iscsi
etc/xen/scripts/block-nbd
etc/xen/scripts/block-tap
etc/xen/scripts/colo-proxy-setup
etc/xen/scripts/external-device-migrate
etc/xen/scripts/hotplugpath.sh
etc/xen/scripts/launch-xenstore
etc/xen/scripts/locking.sh
etc/xen/scripts/logging.sh
etc/xen/scripts/remus-netbuf-setup
etc/xen/scripts/vif-bridge
etc/xen/scripts/vif-common.sh
etc/xen/scripts/vif-nat
etc/xen/scripts/vif-openvswitch
etc/xen/scripts/vif-route
etc/xen/scripts/vif-setup
etc/xen/scripts/vscsi
etc/xen/scripts/xen-hotplug-common.sh
etc/xen/scripts/xen-network-common.sh
etc/xen/scripts/xen-script-common.sh
etc/xen/xl.conf
etc/xen/xlexample.hvm
etc/xen/xlexample.pvhlinux
etc/xen/xlexample.pvlinux
run/
run/xen/
run/xenstored/
usr/
usr/bin/
usr/bin/flask-get-bool
usr/bin/flask-getenforce
usr/bin/flask-label-pci
usr/bin/flask-loadpolicy
usr/bin/flask-set-bool
usr/bin/flask-setenforce
usr/bin/gdbsx
usr/bin/pygrub
usr/bin/vchan-socket-proxy
usr/bin/xen-access
usr/bin/xen-cpuid
usr/bin/xen-detect
usr/bin/xen-diag
usr/bin/xen-hptool
usr/bin/xen-hvmcrash
usr/bin/xen-hvmctx
usr/bin/xen-kdd
usr/bin/xen-livepatch
usr/bin/xen-lowmemd
usr/bin/xen-memshare
usr/bin/xen-mfndump
usr/bin/xen-ucode
usr/bin/xenalyze
usr/bin/xenbaked
usr/bin/xencons
usr/bin/xenconsoled
usr/bin/xencov
usr/bin/xencov_split
usr/bin/xenhypfs
usr/bin/xenlockprof
usr/bin/xenmon
usr/bin/xenperf
usr/bin/xenpm
usr/bin/xenpmd
usr/bin/xenstore
usr/bin/xenstore-chmod
usr/bin/xenstore-control
usr/bin/xenstore-exists
usr/bin/xenstore-list
usr/bin/xenstore-ls
usr/bin/xenstore-read
usr/bin/xenstore-rm
usr/bin/xenstore-watch
usr/bin/xenstore-write
usr/bin/xenstored
usr/bin/xentop
usr/bin/xentrace
usr/bin/xentrace_format
usr/bin/xentrace_setmask
usr/bin/xentrace_setsize
usr/bin/xenwatchdogd
usr/bin/xl
usr/include/
usr/include/_libxl_list.h
usr/include/_libxl_types.h
usr/include/_libxl_types_json.h
usr/include/libxenvchan.h
usr/include/libxl.h
usr/include/libxl_event.h
usr/include/libxl_json.h
usr/include/libxl_utils.h
usr/include/libxl_uuid.h
usr/include/libxlutil.h
usr/include/xen/
usr/include/xen/COPYING
usr/include/xen/arch-arm.h
usr/include/xen/arch-arm/
usr/include/xen/arch-arm/hvm/
usr/include/xen/arch-arm/hvm/save.h
usr/include/xen/arch-arm/smccc.h
usr/include/xen/arch-x86/
usr/include/xen/arch-x86/cpufeatureset.h
usr/include/xen/arch-x86/cpuid.h
usr/include/xen/arch-x86/hvm/
usr/include/xen/arch-x86/hvm/save.h
usr/include/xen/arch-x86/hvm/start_info.h
usr/include/xen/arch-x86/pmu.h
usr/include/xen/arch-x86/xen-mca.h
usr/include/xen/arch-x86/xen-x86_32.h
usr/include/xen/arch-x86/xen-x86_64.h
usr/include/xen/arch-x86/xen.h
usr/include/xen/arch-x86_32.h
usr/include/xen/arch-x86_64.h
usr/include/xen/argo.h
usr/include/xen/callback.h
usr/include/xen/device_tree_defs.h
usr/include/xen/dom0_ops.h
usr/include/xen/domctl.h
usr/include/xen/elfnote.h
usr/include/xen/errno.h
usr/include/xen/event_channel.h
usr/include/xen/features.h
usr/include/xen/foreign/
usr/include/xen/foreign/arm32.h
usr/include/xen/foreign/arm64.h
usr/include/xen/foreign/x86_32.h
usr/include/xen/foreign/x86_64.h
usr/include/xen/grant_table.h
usr/include/xen/hvm/
usr/include/xen/hvm/dm_op.h
usr/include/xen/hvm/e820.h
usr/include/xen/hvm/hvm_info_table.h
usr/include/xen/hvm/hvm_op.h
usr/include/xen/hvm/hvm_vcpu.h
usr/include/xen/hvm/hvm_xs_strings.h
usr/include/xen/hvm/ioreq.h
usr/include/xen/hvm/params.h
usr/include/xen/hvm/pvdrivers.h
usr/include/xen/hvm/save.h
usr/include/xen/hypfs.h
usr/include/xen/io/
usr/include/xen/io/9pfs.h
usr/include/xen/io/blkif.h
usr/include/xen/io/cameraif.h
usr/include/xen/io/console.h
usr/include/xen/io/displif.h
usr/include/xen/io/fbif.h
usr/include/xen/io/fsif.h
usr/include/xen/io/kbdif.h
usr/include/xen/io/libxenvchan.h
usr/include/xen/io/netif.h
usr/include/xen/io/pciif.h
usr/include/xen/io/protocols.h
usr/include/xen/io/pvcalls.h
usr/include/xen/io/ring.h
usr/include/xen/io/sndif.h
usr/include/xen/io/tpmif.h
usr/include/xen/io/usbif.h
usr/include/xen/io/vscsiif.h
usr/include/xen/io/xenbus.h
usr/include/xen/io/xs_wire.h
usr/include/xen/kexec.h
usr/include/xen/memory.h
usr/include/xen/nmi.h
usr/include/xen/physdev.h
usr/include/xen/platform.h
usr/include/xen/pmu.h
usr/include/xen/sched.h
usr/include/xen/sys/
usr/include/xen/sys/evtchn.h
usr/include/xen/sys/gntalloc.h
usr/include/xen/sys/gntdev.h
usr/include/xen/sys/privcmd.h
usr/include/xen/sys/xenbus_dev.h
usr/include/xen/sysctl.h
usr/include/xen/tmem.h
usr/include/xen/trace.h
usr/include/xen/vcpu.h
usr/include/xen/version.h
usr/include/xen/vm_event.h
usr/include/xen/xen-compat.h
usr/include/xen/xen.h
usr/include/xen/xencomm.h
usr/include/xen/xenoprof.h
usr/include/xen/xsm/
usr/include/xen/xsm/flask_op.h
usr/include/xencall.h
usr/include/xenctrl.h
usr/include/xenctrl_compat.h
usr/include/xendevicemodel.h
usr/include/xenevtchn.h
usr/include/xenforeignmemory.h
usr/include/xenfsimage.h
usr/include/xenfsimage_grub.h
usr/include/xenfsimage_plugin.h
usr/include/xengnttab.h
usr/include/xenguest.h
usr/include/xenhypfs.h
usr/include/xenstat.h
usr/include/xenstore-compat/
usr/include/xenstore-compat/xs.h
usr/include/xenstore-compat/xs_lib.h
usr/include/xenstore.h
usr/include/xenstore_lib.h
usr/include/xentoolcore.h
usr/include/xentoolcore_internal.h
usr/include/xentoollog.h
usr/include/xs.h
usr/include/xs_lib.h
usr/lib/
usr/lib/libxencall.so
usr/lib/libxencall.so.1
usr/lib/libxencall.so.1.2
usr/lib/libxenctrl.so
usr/lib/libxenctrl.so.4.15
usr/lib/libxenctrl.so.4.15.0
usr/lib/libxendevicemodel.so
usr/lib/libxendevicemodel.so.1
usr/lib/libxendevicemodel.so.1.4
usr/lib/libxenevtchn.so
usr/lib/libxenevtchn.so.1
usr/lib/libxenevtchn.so.1.2
usr/lib/libxenforeignmemory.so
usr/lib/libxenforeignmemory.so.1
usr/lib/libxenforeignmemory.so.1.4
usr/lib/libxenfsimage.so
usr/lib/libxenfsimage.so.4.15
usr/lib/libxenfsimage.so.4.15.0
usr/lib/libxengnttab.so
usr/lib/libxengnttab.so.1
usr/lib/libxengnttab.so.1.2
usr/lib/libxenguest.so
usr/lib/libxenguest.so.4.15
usr/lib/libxenguest.so.4.15.0
usr/lib/libxenhypfs.so
usr/lib/libxenhypfs.so.1
usr/lib/libxenhypfs.so.1.0
usr/lib/libxenlight.so
usr/lib/libxenlight.so.4.15
usr/lib/libxenlight.so.4.15.0
usr/lib/libxenstat.so
usr/lib/libxenstat.so.4.15
usr/lib/libxenstat.so.4.15.0
usr/lib/libxenstore.so
usr/lib/libxenstore.so.3.0
usr/lib/libxenstore.so.3.0.3
usr/lib/libxentoolcore.so
usr/lib/libxentoolcore.so.1
usr/lib/libxentoolcore.so.1.0
usr/lib/libxentoollog.so
usr/lib/libxentoollog.so.1
usr/lib/libxentoollog.so.1.0
usr/lib/libxenvchan.so
usr/lib/libxenvchan.so.4.15
usr/lib/libxenvchan.so.4.15.0
usr/lib/libxlutil.so
usr/lib/libxlutil.so.4.15
usr/lib/libxlutil.so.4.15.0
usr/lib/modules-load.d/
usr/lib/modules-load.d/xen.conf
usr/lib/pkgconfig/
usr/lib/pkgconfig/xencall.pc
usr/lib/pkgconfig/xencontrol.pc
usr/lib/pkgconfig/xendevicemodel.pc
usr/lib/pkgconfig/xenevtchn.pc
usr/lib/pkgconfig/xenforeignmemory.pc
usr/lib/pkgconfig/xengnttab.pc
usr/lib/pkgconfig/xenguest.pc
usr/lib/pkgconfig/xenhypfs.pc
usr/lib/pkgconfig/xenlight.pc
usr/lib/pkgconfig/xenstat.pc
usr/lib/pkgconfig/xenstore.pc
usr/lib/pkgconfig/xentoolcore.pc
usr/lib/pkgconfig/xentoollog.pc
usr/lib/pkgconfig/xenvchan.pc
usr/lib/pkgconfig/xlutil.pc
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/grub/
usr/lib/python3.9/site-packages/grub/ExtLinuxConf.py
usr/lib/python3.9/site-packages/grub/GrubConf.py
usr/lib/python3.9/site-packages/grub/LiloConf.py
usr/lib/python3.9/site-packages/grub/__init__.py
usr/lib/python3.9/site-packages/grub/__pycache__/
usr/lib/python3.9/site-packages/grub/__pycache__/ExtLinuxConf.cpython-39.pyc
usr/lib/python3.9/site-packages/grub/__pycache__/GrubConf.cpython-39.pyc
usr/lib/python3.9/site-packages/grub/__pycache__/LiloConf.cpython-39.pyc
usr/lib/python3.9/site-packages/grub/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/pygrub-0.6-py3.9.egg-info
usr/lib/python3.9/site-packages/xen-3.0-py3.9.egg-info
usr/lib/python3.9/site-packages/xen/
usr/lib/python3.9/site-packages/xen/__init__.py
usr/lib/python3.9/site-packages/xen/__pycache__/
usr/lib/python3.9/site-packages/xen/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/__pycache__/util.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/lowlevel/
usr/lib/python3.9/site-packages/xen/lowlevel/__init__.py
usr/lib/python3.9/site-packages/xen/lowlevel/__pycache__/
usr/lib/python3.9/site-packages/xen/lowlevel/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/lowlevel/xc.cpython-39-x86_64-linux-gnu.so
usr/lib/python3.9/site-packages/xen/lowlevel/xs.cpython-39-x86_64-linux-gnu.so
usr/lib/python3.9/site-packages/xen/migration/
usr/lib/python3.9/site-packages/xen/migration/__init__.py
usr/lib/python3.9/site-packages/xen/migration/__pycache__/
usr/lib/python3.9/site-packages/xen/migration/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/migration/__pycache__/legacy.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/migration/__pycache__/libxc.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/migration/__pycache__/libxl.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/migration/__pycache__/public.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/migration/__pycache__/tests.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/migration/__pycache__/verify.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/migration/__pycache__/xl.cpython-39.pyc
usr/lib/python3.9/site-packages/xen/migration/legacy.py
usr/lib/python3.9/site-packages/xen/migration/libxc.py
usr/lib/python3.9/site-packages/xen/migration/libxl.py
usr/lib/python3.9/site-packages/xen/migration/public.py
usr/lib/python3.9/site-packages/xen/migration/tests.py
usr/lib/python3.9/site-packages/xen/migration/verify.py
usr/lib/python3.9/site-packages/xen/migration/xl.py
usr/lib/python3.9/site-packages/xen/util.py
usr/lib/python3.9/site-packages/xenfsimage.cpython-39-x86_64-linux-gnu.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/proc-xen.mount
usr/lib/systemd/system/var-lib-xenstored.mount
usr/lib/systemd/system/xen-init-dom0.service
usr/lib/systemd/system/xen-qemu-dom0-disk-backend.service
usr/lib/systemd/system/xen-watchdog.service
usr/lib/systemd/system/xenconsoled.service
usr/lib/systemd/system/xendomains.service
usr/lib/systemd/system/xendriverdomain.service
usr/lib/systemd/system/xenstored.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/xen.conf
usr/lib/xen/
usr/lib/xen/bin/
usr/lib/xen/bin/convert-legacy-stream
usr/lib/xen/bin/depriv-fd-checker
usr/lib/xen/bin/elf2dmp
usr/lib/xen/bin/init-xenstore-domain
usr/lib/xen/bin/ivshmem-client
usr/lib/xen/bin/ivshmem-server
usr/lib/xen/bin/libxl-save-helper
usr/lib/xen/bin/lsevtchn
usr/lib/xen/bin/pygrub
usr/lib/xen/bin/qemu-edid
usr/lib/xen/bin/qemu-img
usr/lib/xen/bin/qemu-io
usr/lib/xen/bin/qemu-nbd
usr/lib/xen/bin/qemu-storage-daemon
usr/lib/xen/bin/qemu-system-i386
usr/lib/xen/bin/readnotes
usr/lib/xen/bin/verify-stream-v2
usr/lib/xen/bin/xen-init-dom0
usr/lib/xen/bin/xenconsole
usr/lib/xen/bin/xenctx
usr/lib/xen/bin/xendomains
usr/lib/xen/bin/xenpaging
usr/lib/xen/bin/xenpvnetboot
usr/lib/xen/boot/
usr/lib/xen/boot/hvmloader
usr/lib/xen/boot/xen-shim
usr/lib/xen/include/
usr/lib/xen/libexec/
usr/lib/xen/libexec/qemu-bridge-helper
usr/lib/xen/libexec/qemu-pr-helper
usr/lib/xen/libexec/virtfs-proxy-helper
usr/lib/xen/libexec/virtiofsd
usr/lib/xenfsimage/
usr/lib/xenfsimage/ext2fs-lib/
usr/lib/xenfsimage/ext2fs-lib/fsimage.so
usr/lib/xenfsimage/fat/
usr/lib/xenfsimage/fat/fsimage.so
usr/lib/xenfsimage/iso9660/
usr/lib/xenfsimage/iso9660/fsimage.so
usr/lib/xenfsimage/reiserfs/
usr/lib/xenfsimage/reiserfs/fsimage.so
usr/lib/xenfsimage/ufs/
usr/lib/xenfsimage/ufs/fsimage.so
usr/lib/xenfsimage/xfs/
usr/lib/xenfsimage/xfs/fsimage.so
usr/lib/xenfsimage/zfs/
usr/lib/xenfsimage/zfs/fsimage.so
usr/share/
usr/share/libalpm/
usr/share/libalpm/hooks/
usr/share/libalpm/hooks/xen-amd-ucode.hook
usr/share/libalpm/hooks/xen-intel-ucode.hook
usr/share/libalpm/scripts/
usr/share/libalpm/scripts/xen-ucode-extract.sh
usr/share/qemu-xen/
usr/share/qemu-xen/applications/
usr/share/qemu-xen/applications/qemu.desktop
usr/share/qemu-xen/icons/
usr/share/qemu-xen/icons/hicolor/
usr/share/qemu-xen/icons/hicolor/128x128/
usr/share/qemu-xen/icons/hicolor/128x128/apps/
usr/share/qemu-xen/icons/hicolor/128x128/apps/qemu.png
usr/share/qemu-xen/icons/hicolor/16x16/
usr/share/qemu-xen/icons/hicolor/16x16/apps/
usr/share/qemu-xen/icons/hicolor/16x16/apps/qemu.png
usr/share/qemu-xen/icons/hicolor/24x24/
usr/share/qemu-xen/icons/hicolor/24x24/apps/
usr/share/qemu-xen/icons/hicolor/24x24/apps/qemu.png
usr/share/qemu-xen/icons/hicolor/256x256/
usr/share/qemu-xen/icons/hicolor/256x256/apps/
usr/share/qemu-xen/icons/hicolor/256x256/apps/qemu.png
usr/share/qemu-xen/icons/hicolor/32x32/
usr/share/qemu-xen/icons/hicolor/32x32/apps/
usr/share/qemu-xen/icons/hicolor/32x32/apps/qemu.bmp
usr/share/qemu-xen/icons/hicolor/32x32/apps/qemu.png
usr/share/qemu-xen/icons/hicolor/48x48/
usr/share/qemu-xen/icons/hicolor/48x48/apps/
usr/share/qemu-xen/icons/hicolor/48x48/apps/qemu.png
usr/share/qemu-xen/icons/hicolor/512x512/
usr/share/qemu-xen/icons/hicolor/512x512/apps/
usr/share/qemu-xen/icons/hicolor/512x512/apps/qemu.png
usr/share/qemu-xen/icons/hicolor/64x64/
usr/share/qemu-xen/icons/hicolor/64x64/apps/
usr/share/qemu-xen/icons/hicolor/64x64/apps/qemu.png
usr/share/qemu-xen/icons/hicolor/scalable/
usr/share/qemu-xen/icons/hicolor/scalable/apps/
usr/share/qemu-xen/icons/hicolor/scalable/apps/qemu.svg
usr/share/qemu-xen/qemu/
usr/share/qemu-xen/qemu/QEMU,cgthree.bin
usr/share/qemu-xen/qemu/QEMU,tcx.bin
usr/share/qemu-xen/qemu/bamboo.dtb
usr/share/qemu-xen/qemu/bios-256k.bin
usr/share/qemu-xen/qemu/bios-microvm.bin
usr/share/qemu-xen/qemu/bios.bin
usr/share/qemu-xen/qemu/canyonlands.dtb
usr/share/qemu-xen/qemu/edk2-aarch64-code.fd
usr/share/qemu-xen/qemu/edk2-arm-code.fd
usr/share/qemu-xen/qemu/edk2-arm-vars.fd
usr/share/qemu-xen/qemu/edk2-i386-code.fd
usr/share/qemu-xen/qemu/edk2-i386-secure-code.fd
usr/share/qemu-xen/qemu/edk2-i386-vars.fd
usr/share/qemu-xen/qemu/edk2-licenses.txt
usr/share/qemu-xen/qemu/edk2-x86_64-code.fd
usr/share/qemu-xen/qemu/edk2-x86_64-secure-code.fd
usr/share/qemu-xen/qemu/efi-e1000.rom
usr/share/qemu-xen/qemu/efi-e1000e.rom
usr/share/qemu-xen/qemu/efi-eepro100.rom
usr/share/qemu-xen/qemu/efi-ne2k_pci.rom
usr/share/qemu-xen/qemu/efi-pcnet.rom
usr/share/qemu-xen/qemu/efi-rtl8139.rom
usr/share/qemu-xen/qemu/efi-virtio.rom
usr/share/qemu-xen/qemu/efi-vmxnet3.rom
usr/share/qemu-xen/qemu/firmware/
usr/share/qemu-xen/qemu/firmware/50-edk2-i386-secure.json
usr/share/qemu-xen/qemu/firmware/50-edk2-x86_64-secure.json
usr/share/qemu-xen/qemu/firmware/60-edk2-aarch64.json
usr/share/qemu-xen/qemu/firmware/60-edk2-arm.json
usr/share/qemu-xen/qemu/firmware/60-edk2-i386.json
usr/share/qemu-xen/qemu/firmware/60-edk2-x86_64.json
usr/share/qemu-xen/qemu/hppa-firmware.img
usr/share/qemu-xen/qemu/keymaps/
usr/share/qemu-xen/qemu/keymaps/ar
usr/share/qemu-xen/qemu/keymaps/bepo
usr/share/qemu-xen/qemu/keymaps/cz
usr/share/qemu-xen/qemu/keymaps/da
usr/share/qemu-xen/qemu/keymaps/de
usr/share/qemu-xen/qemu/keymaps/de-ch
usr/share/qemu-xen/qemu/keymaps/en-gb
usr/share/qemu-xen/qemu/keymaps/en-us
usr/share/qemu-xen/qemu/keymaps/es
usr/share/qemu-xen/qemu/keymaps/et
usr/share/qemu-xen/qemu/keymaps/fi
usr/share/qemu-xen/qemu/keymaps/fo
usr/share/qemu-xen/qemu/keymaps/fr
usr/share/qemu-xen/qemu/keymaps/fr-be
usr/share/qemu-xen/qemu/keymaps/fr-ca
usr/share/qemu-xen/qemu/keymaps/fr-ch
usr/share/qemu-xen/qemu/keymaps/hr
usr/share/qemu-xen/qemu/keymaps/hu
usr/share/qemu-xen/qemu/keymaps/is
usr/share/qemu-xen/qemu/keymaps/it
usr/share/qemu-xen/qemu/keymaps/ja
usr/share/qemu-xen/qemu/keymaps/lt
usr/share/qemu-xen/qemu/keymaps/lv
usr/share/qemu-xen/qemu/keymaps/mk
usr/share/qemu-xen/qemu/keymaps/nl
usr/share/qemu-xen/qemu/keymaps/no
usr/share/qemu-xen/qemu/keymaps/pl
usr/share/qemu-xen/qemu/keymaps/pt
usr/share/qemu-xen/qemu/keymaps/pt-br
usr/share/qemu-xen/qemu/keymaps/ru
usr/share/qemu-xen/qemu/keymaps/sl
usr/share/qemu-xen/qemu/keymaps/sv
usr/share/qemu-xen/qemu/keymaps/th
usr/share/qemu-xen/qemu/keymaps/tr
usr/share/qemu-xen/qemu/kvmvapic.bin
usr/share/qemu-xen/qemu/linuxboot.bin
usr/share/qemu-xen/qemu/linuxboot_dma.bin
usr/share/qemu-xen/qemu/multiboot.bin
usr/share/qemu-xen/qemu/openbios-ppc
usr/share/qemu-xen/qemu/openbios-sparc32
usr/share/qemu-xen/qemu/openbios-sparc64
usr/share/qemu-xen/qemu/opensbi-riscv32-sifive_u-fw_jump.bin
usr/share/qemu-xen/qemu/opensbi-riscv32-virt-fw_jump.bin
usr/share/qemu-xen/qemu/opensbi-riscv64-sifive_u-fw_jump.bin
usr/share/qemu-xen/qemu/opensbi-riscv64-virt-fw_jump.bin
usr/share/qemu-xen/qemu/palcode-clipper
usr/share/qemu-xen/qemu/petalogix-ml605.dtb
usr/share/qemu-xen/qemu/petalogix-s3adsp1800.dtb
usr/share/qemu-xen/qemu/pvh.bin
usr/share/qemu-xen/qemu/pxe-e1000.rom
usr/share/qemu-xen/qemu/pxe-eepro100.rom
usr/share/qemu-xen/qemu/pxe-ne2k_pci.rom
usr/share/qemu-xen/qemu/pxe-pcnet.rom
usr/share/qemu-xen/qemu/pxe-rtl8139.rom
usr/share/qemu-xen/qemu/pxe-virtio.rom
usr/share/qemu-xen/qemu/qemu-nsis.bmp
usr/share/qemu-xen/qemu/qemu_vga.ndrv
usr/share/qemu-xen/qemu/s390-ccw.img
usr/share/qemu-xen/qemu/s390-netboot.img
usr/share/qemu-xen/qemu/sgabios.bin
usr/share/qemu-xen/qemu/skiboot.lid
usr/share/qemu-xen/qemu/slof.bin
usr/share/qemu-xen/qemu/trace-events-all
usr/share/qemu-xen/qemu/u-boot-sam460-20100605.bin
usr/share/qemu-xen/qemu/u-boot.e500
usr/share/qemu-xen/qemu/vgabios-ati.bin
usr/share/qemu-xen/qemu/vgabios-bochs-display.bin
usr/share/qemu-xen/qemu/vgabios-cirrus.bin
usr/share/qemu-xen/qemu/vgabios-qxl.bin
usr/share/qemu-xen/qemu/vgabios-ramfb.bin
usr/share/qemu-xen/qemu/vgabios-stdvga.bin
usr/share/qemu-xen/qemu/vgabios-virtio.bin
usr/share/qemu-xen/qemu/vgabios-vmware.bin
usr/share/qemu-xen/qemu/vgabios.bin
usr/share/qemu-xen/qemu/vhost-user/
usr/share/qemu-xen/qemu/vhost-user/50-qemu-virtiofsd.json
var/
var/lib/
var/lib/xen/
var/lib/xen/dump/
var/lib/xen/xenpaging/
var/lib/xenstored/
var/log/
var/log/xen/
var/log/xen/console/
