# files/rustdesk-1.1.9-2/files
%FILES%
usr/
usr/bin/
usr/bin/rustdesk
usr/lib/
usr/lib/rustdesk/
usr/lib/rustdesk/libsciter-gtk.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rustdesk.service
usr/share/
usr/share/applications/
usr/share/applications/rustdesk.desktop
usr/share/licenses/
usr/share/licenses/rustdesk/
usr/share/licenses/rustdesk/LICENSE
usr/share/pixmaps/
usr/share/pixmaps/rustdesk.png
