# files/key-mapper-git-1.2.1.r14.g3fb42ad-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/keymapper.Control.conf
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/key-mapper-autoload.desktop
usr/
usr/bin/
usr/bin/key-mapper-control
usr/bin/key-mapper-gtk
usr/bin/key-mapper-helper
usr/bin/key-mapper-service
usr/lib/
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/key_mapper-1.2.1-py3.9.egg-info/
usr/lib/python3.9/site-packages/key_mapper-1.2.1-py3.9.egg-info/PKG-INFO
usr/lib/python3.9/site-packages/key_mapper-1.2.1-py3.9.egg-info/SOURCES.txt
usr/lib/python3.9/site-packages/key_mapper-1.2.1-py3.9.egg-info/dependency_links.txt
usr/lib/python3.9/site-packages/key_mapper-1.2.1-py3.9.egg-info/requires.txt
usr/lib/python3.9/site-packages/key_mapper-1.2.1-py3.9.egg-info/top_level.txt
usr/lib/python3.9/site-packages/keymapper/
usr/lib/python3.9/site-packages/keymapper/__init__.py
usr/lib/python3.9/site-packages/keymapper/__pycache__/
usr/lib/python3.9/site-packages/keymapper/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/commit_hash.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/config.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/daemon.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/data.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/groups.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/key.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/logger.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/mapping.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/paths.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/presets.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/system_mapping.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/user.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/__pycache__/utils.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/commit_hash.py
usr/lib/python3.9/site-packages/keymapper/config.py
usr/lib/python3.9/site-packages/keymapper/daemon.py
usr/lib/python3.9/site-packages/keymapper/data.py
usr/lib/python3.9/site-packages/keymapper/groups.py
usr/lib/python3.9/site-packages/keymapper/gui/
usr/lib/python3.9/site-packages/keymapper/gui/__init__.py
usr/lib/python3.9/site-packages/keymapper/gui/__pycache__/
usr/lib/python3.9/site-packages/keymapper/gui/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/gui/__pycache__/custom_mapping.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/gui/__pycache__/helper.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/gui/__pycache__/reader.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/gui/__pycache__/row.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/gui/__pycache__/window.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/gui/custom_mapping.py
usr/lib/python3.9/site-packages/keymapper/gui/helper.py
usr/lib/python3.9/site-packages/keymapper/gui/reader.py
usr/lib/python3.9/site-packages/keymapper/gui/row.py
usr/lib/python3.9/site-packages/keymapper/gui/window.py
usr/lib/python3.9/site-packages/keymapper/injection/
usr/lib/python3.9/site-packages/keymapper/injection/__init__.py
usr/lib/python3.9/site-packages/keymapper/injection/__pycache__/
usr/lib/python3.9/site-packages/keymapper/injection/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/__pycache__/consumer_control.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/__pycache__/context.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/__pycache__/injector.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/__pycache__/numlock.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/consumer_control.py
usr/lib/python3.9/site-packages/keymapper/injection/consumers/
usr/lib/python3.9/site-packages/keymapper/injection/consumers/__init__.py
usr/lib/python3.9/site-packages/keymapper/injection/consumers/__pycache__/
usr/lib/python3.9/site-packages/keymapper/injection/consumers/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/consumers/__pycache__/consumer.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/consumers/__pycache__/joystick_to_mouse.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/consumers/__pycache__/keycode_mapper.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/consumers/consumer.py
usr/lib/python3.9/site-packages/keymapper/injection/consumers/joystick_to_mouse.py
usr/lib/python3.9/site-packages/keymapper/injection/consumers/keycode_mapper.py
usr/lib/python3.9/site-packages/keymapper/injection/context.py
usr/lib/python3.9/site-packages/keymapper/injection/injector.py
usr/lib/python3.9/site-packages/keymapper/injection/macros/
usr/lib/python3.9/site-packages/keymapper/injection/macros/__init__.py
usr/lib/python3.9/site-packages/keymapper/injection/macros/__pycache__/
usr/lib/python3.9/site-packages/keymapper/injection/macros/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/macros/__pycache__/macro.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/macros/__pycache__/parse.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/injection/macros/macro.py
usr/lib/python3.9/site-packages/keymapper/injection/macros/parse.py
usr/lib/python3.9/site-packages/keymapper/injection/numlock.py
usr/lib/python3.9/site-packages/keymapper/ipc/
usr/lib/python3.9/site-packages/keymapper/ipc/__init__.py
usr/lib/python3.9/site-packages/keymapper/ipc/__pycache__/
usr/lib/python3.9/site-packages/keymapper/ipc/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/ipc/__pycache__/pipe.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/ipc/__pycache__/shared_dict.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/ipc/__pycache__/socket.cpython-39.pyc
usr/lib/python3.9/site-packages/keymapper/ipc/pipe.py
usr/lib/python3.9/site-packages/keymapper/ipc/shared_dict.py
usr/lib/python3.9/site-packages/keymapper/ipc/socket.py
usr/lib/python3.9/site-packages/keymapper/key.py
usr/lib/python3.9/site-packages/keymapper/logger.py
usr/lib/python3.9/site-packages/keymapper/mapping.py
usr/lib/python3.9/site-packages/keymapper/paths.py
usr/lib/python3.9/site-packages/keymapper/presets.py
usr/lib/python3.9/site-packages/keymapper/system_mapping.py
usr/lib/python3.9/site-packages/keymapper/user.py
usr/lib/python3.9/site-packages/keymapper/utils.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/key-mapper.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/99-key-mapper.rules
usr/share/
usr/share/applications/
usr/share/applications/key-mapper.desktop
usr/share/key-mapper/
usr/share/key-mapper/99-key-mapper.rules
usr/share/key-mapper/key-mapper-autoload.desktop
usr/share/key-mapper/key-mapper-large.png
usr/share/key-mapper/key-mapper.desktop
usr/share/key-mapper/key-mapper.glade
usr/share/key-mapper/key-mapper.policy
usr/share/key-mapper/key-mapper.service
usr/share/key-mapper/key-mapper.svg
usr/share/key-mapper/keymapper.Control.conf
usr/share/key-mapper/lang/
usr/share/key-mapper/lang/it_IT/
usr/share/key-mapper/lang/it_IT/LC_MESSAGES/
usr/share/key-mapper/lang/it_IT/LC_MESSAGES/key-mapper.mo
usr/share/key-mapper/lang/sk_SK/
usr/share/key-mapper/lang/sk_SK/LC_MESSAGES/
usr/share/key-mapper/lang/sk_SK/LC_MESSAGES/key-mapper.mo
usr/share/key-mapper/style.css
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/key-mapper.policy
