# files/php80-fpm-8.0.22-1/files
%FILES%
etc/
etc/php80/
etc/php80/php-fpm.conf
etc/php80/php-fpm.d/
etc/php80/php-fpm.d/www.conf
usr/
usr/bin/
usr/bin/php-fpm80
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/php80-fpm.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/php80-fpm.conf
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/php-fpm80.8.gz
usr/share/php80/
usr/share/php80/fpm/
usr/share/php80/fpm/status.html
