# files/openswan-3.0.0-1/files
%FILES%
etc/
etc/ipsec.conf
etc/ipsec.d/
etc/ipsec.d/aacerts/
etc/ipsec.d/cacerts/
etc/ipsec.d/certs/
etc/ipsec.d/crls/
etc/ipsec.d/examples/
etc/ipsec.d/examples/hub-spoke.conf
etc/ipsec.d/examples/ipv6.conf
etc/ipsec.d/examples/l2tp-cert.conf
etc/ipsec.d/examples/l2tp-psk.conf
etc/ipsec.d/examples/linux-linux.conf
etc/ipsec.d/examples/mast-l2tp-psk.conf
etc/ipsec.d/examples/oe-exclude-dns.conf
etc/ipsec.d/examples/sysctl.conf
etc/ipsec.d/examples/xauth.conf
etc/ipsec.d/ocspcerts/
etc/ipsec.d/policies/
etc/ipsec.d/policies/block
etc/ipsec.d/policies/clear
etc/ipsec.d/policies/clear-or-private
etc/ipsec.d/policies/private
etc/ipsec.d/policies/private-or-clear
etc/ipsec.d/private/
usr/
usr/bin/
usr/bin/ipsec
usr/lib/
usr/lib/ipsec/
usr/lib/ipsec/_copyright
usr/lib/ipsec/_include
usr/lib/ipsec/_keycensor
usr/lib/ipsec/_plutoload
usr/lib/ipsec/_plutorun
usr/lib/ipsec/_realsetup
usr/lib/ipsec/_secretcensor
usr/lib/ipsec/_startklips
usr/lib/ipsec/_startnetkey
usr/lib/ipsec/_updown
usr/lib/ipsec/_updown.klips
usr/lib/ipsec/_updown.mast
usr/lib/ipsec/_updown.netkey
usr/lib/openswan/
usr/lib/openswan/addconn
usr/lib/openswan/auto
usr/lib/openswan/barf
usr/lib/openswan/eroute
usr/lib/openswan/ikeping
usr/lib/openswan/initnss
usr/lib/openswan/klipsdebug
usr/lib/openswan/look
usr/lib/openswan/newhostkey
usr/lib/openswan/pf_key
usr/lib/openswan/pluto
usr/lib/openswan/policy
usr/lib/openswan/ranbits
usr/lib/openswan/rsasigkey
usr/lib/openswan/secrets
usr/lib/openswan/setup
usr/lib/openswan/showdefaults
usr/lib/openswan/showhostkey
usr/lib/openswan/spi
usr/lib/openswan/spigrp
usr/lib/openswan/status
usr/lib/openswan/tncfg
usr/lib/openswan/verify
usr/lib/openswan/whack
usr/lib/systemd/
usr/lib/systemd/scripts/
usr/lib/systemd/scripts/ipsec
usr/lib/systemd/system/
usr/lib/systemd/system/openswan.service
usr/share/
usr/share/doc/
usr/share/doc/openswan/
usr/share/doc/openswan/ipsec.conf-sample
usr/share/licenses/
usr/share/licenses/openswan/
usr/share/licenses/openswan/LICENSE
usr/share/man/
usr/share/man/man3/
usr/share/man/man5/
usr/share/man/man5/ipsec.conf.5.gz
usr/share/man/man5/ipsec.secrets.5.gz
usr/share/man/man5/ipsec_eroute.5.gz
usr/share/man/man5/ipsec_klipsdebug.5.gz
usr/share/man/man5/ipsec_pf_key.5.gz
usr/share/man/man5/ipsec_spi.5.gz
usr/share/man/man5/ipsec_spigrp.5.gz
usr/share/man/man5/ipsec_tncfg.5.gz
usr/share/man/man5/ipsec_trap_count.5.gz
usr/share/man/man5/ipsec_trap_sendcount.5.gz
usr/share/man/man5/ipsec_version.5.gz
usr/share/man/man8/
usr/share/man/man8/ipsec.8.gz
usr/share/man/man8/ipsec__copyright.8.gz
usr/share/man/man8/ipsec__include.8.gz
usr/share/man/man8/ipsec__keycensor.8.gz
usr/share/man/man8/ipsec__plutoload.8.gz
usr/share/man/man8/ipsec__plutorun.8.gz
usr/share/man/man8/ipsec__realsetup.8.gz
usr/share/man/man8/ipsec__secretcensor.8.gz
usr/share/man/man8/ipsec__startklips.8.gz
usr/share/man/man8/ipsec__startnetkey.8.gz
usr/share/man/man8/ipsec__updown.8.gz
usr/share/man/man8/ipsec__updown.klips.8.gz
usr/share/man/man8/ipsec__updown.mast.8.gz
usr/share/man/man8/ipsec__updown.netkey.8.gz
usr/share/man/man8/ipsec_addconn.8.gz
usr/share/man/man8/ipsec_auto.8.gz
usr/share/man/man8/ipsec_barf.8.gz
usr/share/man/man8/ipsec_eroute.8.gz
usr/share/man/man8/ipsec_ikeping.8.gz
usr/share/man/man8/ipsec_initnss.8.gz
usr/share/man/man8/ipsec_klipsdebug.8.gz
usr/share/man/man8/ipsec_look.8.gz
usr/share/man/man8/ipsec_newhostkey.8.gz
usr/share/man/man8/ipsec_pf_key.8.gz
usr/share/man/man8/ipsec_pluto.8.gz
usr/share/man/man8/ipsec_policy.8.gz
usr/share/man/man8/ipsec_ranbits.8.gz
usr/share/man/man8/ipsec_rsasigkey.8.gz
usr/share/man/man8/ipsec_secrets.8.gz
usr/share/man/man8/ipsec_setup.8.gz
usr/share/man/man8/ipsec_showdefaults.8.gz
usr/share/man/man8/ipsec_showhostkey.8.gz
usr/share/man/man8/ipsec_spi.8.gz
usr/share/man/man8/ipsec_spigrp.8.gz
usr/share/man/man8/ipsec_status.8.gz
usr/share/man/man8/ipsec_tncfg.8.gz
usr/share/man/man8/ipsec_verify.8.gz
usr/share/man/man8/ipsec_whack.8.gz
