# files/polyglot-3.5-1/files
%FILES%
usr/
usr/bin/
usr/bin/polyglot
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/polyglot.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/polyglot.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/polyglot.conf
