# files/appmenu-gtk-module-git-0.7.6.r21.g87abe53-1/files
%FILES%
etc/
etc/X11/
etc/X11/xinit/
etc/X11/xinit/xinitrc.d/
etc/X11/xinit/xinitrc.d/80appmenu-gtk-module
usr/
usr/include/
usr/include/appmenu-gtk-parser/
usr/include/appmenu-gtk-parser/appmenu-gtk-action-group.h
usr/include/appmenu-gtk-parser/appmenu-gtk-menu-shell.h
usr/include/appmenu-gtk-parser/appmenu-gtk-parser.h
usr/lib/
usr/lib/gtk-2.0/
usr/lib/gtk-2.0/modules/
usr/lib/gtk-2.0/modules/libappmenu-gtk-module.so
usr/lib/gtk-3.0/
usr/lib/gtk-3.0/modules/
usr/lib/gtk-3.0/modules/libappmenu-gtk-module.so
usr/lib/libappmenu-gtk2-parser.so
usr/lib/libappmenu-gtk2-parser.so.0
usr/lib/libappmenu-gtk2-parser.so.0.7
usr/lib/libappmenu-gtk3-parser.so
usr/lib/libappmenu-gtk3-parser.so.0
usr/lib/libappmenu-gtk3-parser.so.0.7
usr/lib/pkgconfig/
usr/lib/pkgconfig/appmenu-gtk2-parser.pc
usr/lib/pkgconfig/appmenu-gtk3-parser.pc
usr/lib/systemd/
usr/lib/systemd/user/
usr/lib/systemd/user/appmenu-gtk-module.service
usr/share/
usr/share/appmenu-gtk-module/
usr/share/appmenu-gtk-module/doc/
usr/share/appmenu-gtk-module/doc/AUTHORS
usr/share/appmenu-gtk-module/doc/LICENSE
usr/share/appmenu-gtk-module/doc/README.md
usr/share/appmenu-gtk-module/doc/TODO
usr/share/glib-2.0/
usr/share/glib-2.0/schemas/
usr/share/glib-2.0/schemas/org.appmenu.gtk-module.gschema.xml
