# files/mullvad-vpn-cli-2022.1-1/files
%FILES%
etc/
etc/systemd/
etc/systemd/system/
etc/systemd/system/mullvad-daemon.service.d/
etc/systemd/system/mullvad-daemon.service.d/override.conf
opt/
opt/mullvad-vpn-cli/
opt/mullvad-vpn-cli/api-ip-address.txt
opt/mullvad-vpn-cli/binaries/
opt/mullvad-vpn-cli/binaries/.git
opt/mullvad-vpn-cli/binaries/.gitattributes
opt/mullvad-vpn-cli/binaries/.gitignore
opt/mullvad-vpn-cli/binaries/.gitmodules
opt/mullvad-vpn-cli/binaries/Makefile
opt/mullvad-vpn-cli/binaries/README.md
opt/mullvad-vpn-cli/binaries/aarch64-apple-darwin/
opt/mullvad-vpn-cli/binaries/aarch64-apple-darwin/openvpn
opt/mullvad-vpn-cli/binaries/aarch64-apple-darwin/sslocal
opt/mullvad-vpn-cli/binaries/libmnl/
opt/mullvad-vpn-cli/binaries/libnftnl/
opt/mullvad-vpn-cli/binaries/libsodium/
opt/mullvad-vpn-cli/binaries/lz4/
opt/mullvad-vpn-cli/binaries/openssl/
opt/mullvad-vpn-cli/binaries/openvpn-build/
opt/mullvad-vpn-cli/binaries/openvpn/
opt/mullvad-vpn-cli/binaries/setup-generic-buildsystem.6.sh
opt/mullvad-vpn-cli/binaries/shadowsocks-rust/
opt/mullvad-vpn-cli/binaries/shadowsocks-windows.bat
opt/mullvad-vpn-cli/binaries/win-shadowsocks-cargo-config
opt/mullvad-vpn-cli/binaries/win-split-tunnel/
opt/mullvad-vpn-cli/binaries/wireguard-nt/
opt/mullvad-vpn-cli/binaries/wireguard/
opt/mullvad-vpn-cli/binaries/wireguard/download-wintun.bat
opt/mullvad-vpn-cli/binaries/wireguard/download-wireguard-nt.bat
opt/mullvad-vpn-cli/binaries/x86_64-apple-darwin/
opt/mullvad-vpn-cli/binaries/x86_64-apple-darwin/libsodium.a
opt/mullvad-vpn-cli/binaries/x86_64-apple-darwin/openvpn
opt/mullvad-vpn-cli/binaries/x86_64-apple-darwin/sslocal
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/libsodium.lib
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/nsis/
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/nsis/api.h
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/nsis/nsis_tchar.h
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/nsis/pluginapi-x86-ansi.lib
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/nsis/pluginapi-x86-unicode.lib
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/nsis/pluginapi.h
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/openvpn.exe
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/split-tunnel/
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/split-tunnel/mullvad-split-tunnel.cat
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/split-tunnel/mullvad-split-tunnel.inf
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/split-tunnel/mullvad-split-tunnel.pdb
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/split-tunnel/mullvad-split-tunnel.sys
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/sslocal.exe
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/wintun/
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/wintun/wintun.dll
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/wintun/wintun.h
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/wireguard-nt/
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/wireguard-nt/mullvad-wireguard.dll
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/wireguard-nt/mullvad-wireguard.pdb
opt/mullvad-vpn-cli/binaries/x86_64-pc-windows-msvc/wireguard-nt/wireguard.h
opt/mullvad-vpn-cli/binaries/x86_64-unknown-linux-gnu/
opt/mullvad-vpn-cli/binaries/x86_64-unknown-linux-gnu/libmnl.a
opt/mullvad-vpn-cli/binaries/x86_64-unknown-linux-gnu/libnftnl.a
opt/mullvad-vpn-cli/binaries/x86_64-unknown-linux-gnu/libsodium.a
opt/mullvad-vpn-cli/binaries/x86_64-unknown-linux-gnu/openvpn
opt/mullvad-vpn-cli/binaries/x86_64-unknown-linux-gnu/sslocal
opt/mullvad-vpn-cli/ca.crt
opt/mullvad-vpn-cli/icon-macos.icns
opt/mullvad-vpn-cli/icon.icns
opt/mullvad-vpn-cli/icon.ico
opt/mullvad-vpn-cli/libtalpid_openvpn_plugin.so
opt/mullvad-vpn-cli/linux/
opt/mullvad-vpn-cli/linux/after-install.sh
opt/mullvad-vpn-cli/linux/after-remove.sh
opt/mullvad-vpn-cli/linux/before-install.sh
opt/mullvad-vpn-cli/linux/before-remove.sh
opt/mullvad-vpn-cli/linux/mullvad-daemon.conf
opt/mullvad-vpn-cli/linux/mullvad-daemon.service
opt/mullvad-vpn-cli/linux/mullvad-gui-launcher.sh
opt/mullvad-vpn-cli/linux/post-transaction.sh
opt/mullvad-vpn-cli/linux/problem-report-link
opt/mullvad-vpn-cli/mullvad
opt/mullvad-vpn-cli/mullvad-daemon
opt/mullvad-vpn-cli/mullvad-exclude
opt/mullvad-vpn-cli/mullvad-problem-report
opt/mullvad-vpn-cli/mullvad-setup
opt/mullvad-vpn-cli/pkg-scripts/
opt/mullvad-vpn-cli/pkg-scripts/postinstall
opt/mullvad-vpn-cli/pkg-scripts/preinstall
opt/mullvad-vpn-cli/relays.json
opt/mullvad-vpn-cli/shell-completions/
opt/mullvad-vpn-cli/shell-completions/_mullvad
opt/mullvad-vpn-cli/shell-completions/mullvad.bash
opt/mullvad-vpn-cli/shell-completions/mullvad.fish
opt/mullvad-vpn-cli/uninstall_macos.sh
opt/mullvad-vpn-cli/windows/
opt/mullvad-vpn-cli/windows/installer.nsh
opt/mullvad-vpn-cli/windows/installersidebar.bmp
opt/mullvad-vpn-cli/windows/version.h
usr/
usr/bin/
usr/bin/mullvad
usr/bin/mullvad-exclude
usr/bin/mullvad-problem-report
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mullvad-daemon.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/mullvad
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/mullvad.fish
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_mullvad
