# files/cloudflare-warp-bin-2021.7.2-1/files
%FILES%
usr/
usr/bin/
usr/bin/warp-cli
usr/bin/warp-diag
usr/bin/warp-svc
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/warp-svc.service
usr/usr/
usr/usr/share/
usr/usr/share/doc/
usr/usr/share/doc/cloudflare-warp/
usr/usr/share/doc/cloudflare-warp/changelog.Debian.gz
usr/usr/share/doc/cloudflare-warp/changelog.gz
