# files/macfanctld-0.6-2/files
%FILES%
etc/
etc/macfanctl.conf
usr/
usr/bin/
usr/bin/macfanctld
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/macfanctld.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/macfanctld.1.gz
