# files/nm-iwd-cloud-setup-1.38.2-1/files
%FILES%
usr/
usr/lib/
usr/lib/NetworkManager/
usr/lib/NetworkManager/dispatcher.d/
usr/lib/NetworkManager/dispatcher.d/90-nm-cloud-setup.sh
usr/lib/NetworkManager/dispatcher.d/no-wait.d/
usr/lib/NetworkManager/dispatcher.d/no-wait.d/90-nm-cloud-setup.sh
usr/lib/nm-cloud-setup
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nm-cloud-setup.service
usr/lib/systemd/system/nm-cloud-setup.timer
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/nm-cloud-setup.8.gz
