# files/pacman-auto-update-1.5.3-1/files
%FILES%
usr/
usr/lib/
usr/lib/pacman-auto-update/
usr/lib/pacman-auto-update/pacman-auto-update
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pacman-auto-update.service
usr/lib/systemd/system/pacman-auto-update.timer
