# files/python-validity-0.12-3/files
%FILES%
usr/
usr/bin/
usr/bin/validity-led-dance
usr/bin/validity-sensors-firmware
usr/lib/
usr/lib/python-validity/
usr/lib/python-validity/dbus-service
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/python_validity-0.12-py3.9.egg-info/
usr/lib/python3.9/site-packages/python_validity-0.12-py3.9.egg-info/PKG-INFO
usr/lib/python3.9/site-packages/python_validity-0.12-py3.9.egg-info/SOURCES.txt
usr/lib/python3.9/site-packages/python_validity-0.12-py3.9.egg-info/dependency_links.txt
usr/lib/python3.9/site-packages/python_validity-0.12-py3.9.egg-info/requires.txt
usr/lib/python3.9/site-packages/python_validity-0.12-py3.9.egg-info/top_level.txt
usr/lib/python3.9/site-packages/validitysensor/
usr/lib/python3.9/site-packages/validitysensor/__init__.py
usr/lib/python3.9/site-packages/validitysensor/__pycache__/
usr/lib/python3.9/site-packages/validitysensor/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/blobs.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/blobs_90.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/blobs_97.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/blobs_9a.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/db.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/firmware_tables.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/flash.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/generated_tables.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/hw_tables.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/init.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/init_db.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/init_flash.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/sensor.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/sid.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/table_types.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/timeslot.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/tls.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/upload_fwext.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/usb.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/util.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/__pycache__/winbio_constants.cpython-39.pyc
usr/lib/python3.9/site-packages/validitysensor/blobs.py
usr/lib/python3.9/site-packages/validitysensor/blobs_90.py
usr/lib/python3.9/site-packages/validitysensor/blobs_97.py
usr/lib/python3.9/site-packages/validitysensor/blobs_9a.py
usr/lib/python3.9/site-packages/validitysensor/db.py
usr/lib/python3.9/site-packages/validitysensor/firmware_tables.py
usr/lib/python3.9/site-packages/validitysensor/flash.py
usr/lib/python3.9/site-packages/validitysensor/generated_tables.py
usr/lib/python3.9/site-packages/validitysensor/hw_tables.py
usr/lib/python3.9/site-packages/validitysensor/init.py
usr/lib/python3.9/site-packages/validitysensor/init_db.py
usr/lib/python3.9/site-packages/validitysensor/init_flash.py
usr/lib/python3.9/site-packages/validitysensor/sensor.py
usr/lib/python3.9/site-packages/validitysensor/sid.py
usr/lib/python3.9/site-packages/validitysensor/table_types.py
usr/lib/python3.9/site-packages/validitysensor/timeslot.py
usr/lib/python3.9/site-packages/validitysensor/tls.py
usr/lib/python3.9/site-packages/validitysensor/upload_fwext.py
usr/lib/python3.9/site-packages/validitysensor/usb.py
usr/lib/python3.9/site-packages/validitysensor/util.py
usr/lib/python3.9/site-packages/validitysensor/winbio_constants.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/python3-validity.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-python-validity.rules
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/io.github.uunicorn.Fprint.conf
usr/share/licenses/
usr/share/licenses/python-validity/
usr/share/licenses/python-validity/LICENSE
usr/share/python-validity/
usr/share/python-validity/playground/
usr/share/python-validity/playground/dbus-cmd.py
usr/share/python-validity/playground/factory-reset.py
usr/share/python-validity/playground/lsdbus.py
usr/share/python-validity/playground/prototype.py
