# files/amdgpu-pro-libgl-22.10_1395274-1/files
%FILES%
etc/
etc/amd/
etc/amd/amdapfxx.blb
etc/amd/amdrc
etc/ld.so.conf.d/
etc/ld.so.conf.d/10-amdgpu-pro-x86_64.conf
usr/
usr/bin/
usr/bin/progl
usr/lib/
usr/lib/amdgpu-pro/
usr/lib/amdgpu-pro/libEGL.so
usr/lib/amdgpu-pro/libEGL.so.1
usr/lib/amdgpu-pro/libGL.so
usr/lib/amdgpu-pro/libGL.so.1
usr/lib/amdgpu-pro/libGL.so.1.2
usr/lib/amdgpu-pro/libGLESv2.so
usr/lib/amdgpu-pro/libGLESv2.so.2
usr/lib/amdgpu-pro/libGLX_amd.so
usr/lib/amdgpu-pro/libGLX_amd.so.0
usr/lib/amdgpu-pro/libGLX_amd.so.0.0
usr/lib/amdgpu-pro/libglapi.so
usr/lib/amdgpu-pro/libglapi.so.1
usr/lib/amdgpu-pro/xorg/
usr/lib/amdgpu-pro/xorg/modules/
usr/lib/amdgpu-pro/xorg/modules/extensions/
usr/lib/amdgpu-pro/xorg/modules/extensions/libglx.so
usr/lib/dri/
usr/lib/dri/amdgpu_dri.so
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/progl
usr/share/drirc.d/
usr/share/drirc.d/10-amdgpu-pro.conf.hide
usr/share/licenses/
usr/share/licenses/amdgpu-pro-libgl/
usr/share/licenses/amdgpu-pro-libgl/copyright
