# files/bluez-git-1:5.62.r61.g6ea642f6e-1/files
%FILES%
etc/
etc/bluetooth/
etc/bluetooth/main.conf
usr/
usr/bin/
usr/bin/advtest
usr/bin/amptest
usr/bin/avinfo
usr/bin/avtest
usr/bin/bcmfw
usr/bin/bdaddr
usr/bin/bluemoon
usr/bin/bluetooth-player
usr/bin/bluetoothctl
usr/bin/bneptest
usr/bin/btattach
usr/bin/btconfig
usr/bin/btgatt-client
usr/bin/btgatt-server
usr/bin/btinfo
usr/bin/btiotest
usr/bin/btmgmt
usr/bin/btmon
usr/bin/btpclient
usr/bin/btpclientctl
usr/bin/btproxy
usr/bin/btsnoop
usr/bin/check-selftest
usr/bin/cltest
usr/bin/create-image
usr/bin/eddystone
usr/bin/gatt-service
usr/bin/hcieventmask
usr/bin/hcisecfilter
usr/bin/hex2hcd
usr/bin/hid2hci
usr/bin/hwdb
usr/bin/ibeacon
usr/bin/l2ping
usr/bin/l2test
usr/bin/mcaptest
usr/bin/mpris-proxy
usr/bin/nokfw
usr/bin/obex-client-tool
usr/bin/obex-server-tool
usr/bin/obexctl
usr/bin/oobtest
usr/bin/parse_companies.pl
usr/bin/rctest
usr/bin/rtlfw
usr/bin/scotest
usr/bin/seq2bseq
usr/bin/test-runner
usr/bin/update_compids.sh
usr/include/
usr/include/bluetooth/
usr/include/bluetooth/bluetooth.h
usr/include/bluetooth/bnep.h
usr/include/bluetooth/cmtp.h
usr/include/bluetooth/hci.h
usr/include/bluetooth/hci_lib.h
usr/include/bluetooth/hidp.h
usr/include/bluetooth/l2cap.h
usr/include/bluetooth/rfcomm.h
usr/include/bluetooth/sco.h
usr/include/bluetooth/sdp.h
usr/include/bluetooth/sdp_lib.h
usr/lib/
usr/lib/bluetooth/
usr/lib/bluetooth/bluetoothd
usr/lib/bluetooth/obexd
usr/lib/bluetooth/plugins/
usr/lib/bluetooth/plugins/sixaxis.so
usr/lib/cups/
usr/lib/cups/backend/
usr/lib/cups/backend/bluetooth
usr/lib/libbluetooth.so
usr/lib/libbluetooth.so.3
usr/lib/libbluetooth.so.3.19.6
usr/lib/modprobe.d/
usr/lib/modprobe.d/bluetooth-usb.conf
usr/lib/pkgconfig/
usr/lib/pkgconfig/bluez.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/bluetooth.service
usr/lib/systemd/user/
usr/lib/systemd/user/obex.service
usr/lib/udev/
usr/lib/udev/hid2hci
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/97-hid2hci.rules
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/services/
usr/share/dbus-1/services/org.bluez.obex.service
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.bluez.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/bluetooth.conf
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/btattach.1.gz
usr/share/man/man1/btmon.1.gz
usr/share/man/man1/hid2hci.1.gz
usr/share/man/man1/l2ping.1.gz
usr/share/man/man1/rctest.1.gz
usr/share/man/man8/
usr/share/man/man8/bluetoothd.8.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_bluetoothctl
