# files/ungoogled-chromium-101.0.4951.64-1/files
%FILES%
usr/
usr/bin/
usr/bin/chromedriver
usr/bin/chromium
usr/lib/
usr/lib/chromium/
usr/lib/chromium/chrome-sandbox
usr/lib/chromium/chrome_100_percent.pak
usr/lib/chromium/chrome_200_percent.pak
usr/lib/chromium/chrome_crashpad_handler
usr/lib/chromium/chromium
usr/lib/chromium/libEGL.so
usr/lib/chromium/libGLESv2.so
usr/lib/chromium/libvk_swiftshader.so
usr/lib/chromium/locales/
usr/lib/chromium/locales/af.pak
usr/lib/chromium/locales/am.pak
usr/lib/chromium/locales/ar.pak
usr/lib/chromium/locales/bg.pak
usr/lib/chromium/locales/bn.pak
usr/lib/chromium/locales/ca.pak
usr/lib/chromium/locales/cs.pak
usr/lib/chromium/locales/da.pak
usr/lib/chromium/locales/de.pak
usr/lib/chromium/locales/el.pak
usr/lib/chromium/locales/en-GB.pak
usr/lib/chromium/locales/en-US.pak
usr/lib/chromium/locales/es-419.pak
usr/lib/chromium/locales/es.pak
usr/lib/chromium/locales/et.pak
usr/lib/chromium/locales/fa.pak
usr/lib/chromium/locales/fi.pak
usr/lib/chromium/locales/fil.pak
usr/lib/chromium/locales/fr.pak
usr/lib/chromium/locales/gu.pak
usr/lib/chromium/locales/he.pak
usr/lib/chromium/locales/hi.pak
usr/lib/chromium/locales/hr.pak
usr/lib/chromium/locales/hu.pak
usr/lib/chromium/locales/id.pak
usr/lib/chromium/locales/it.pak
usr/lib/chromium/locales/ja.pak
usr/lib/chromium/locales/kn.pak
usr/lib/chromium/locales/ko.pak
usr/lib/chromium/locales/lt.pak
usr/lib/chromium/locales/lv.pak
usr/lib/chromium/locales/ml.pak
usr/lib/chromium/locales/mr.pak
usr/lib/chromium/locales/ms.pak
usr/lib/chromium/locales/nb.pak
usr/lib/chromium/locales/nl.pak
usr/lib/chromium/locales/pl.pak
usr/lib/chromium/locales/pt-BR.pak
usr/lib/chromium/locales/pt-PT.pak
usr/lib/chromium/locales/ro.pak
usr/lib/chromium/locales/ru.pak
usr/lib/chromium/locales/sk.pak
usr/lib/chromium/locales/sl.pak
usr/lib/chromium/locales/sr.pak
usr/lib/chromium/locales/sv.pak
usr/lib/chromium/locales/sw.pak
usr/lib/chromium/locales/ta.pak
usr/lib/chromium/locales/te.pak
usr/lib/chromium/locales/th.pak
usr/lib/chromium/locales/tr.pak
usr/lib/chromium/locales/uk.pak
usr/lib/chromium/locales/ur.pak
usr/lib/chromium/locales/vi.pak
usr/lib/chromium/locales/zh-CN.pak
usr/lib/chromium/locales/zh-TW.pak
usr/lib/chromium/resources.pak
usr/lib/chromium/swiftshader/
usr/lib/chromium/swiftshader/libEGL.so
usr/lib/chromium/swiftshader/libGLESv2.so
usr/lib/chromium/v8_context_snapshot.bin
usr/lib/chromium/vk_swiftshader_icd.json
usr/share/
usr/share/applications/
usr/share/applications/chromium.desktop
usr/share/drirc.d/
usr/share/drirc.d/10-ungoogled-chromium.conf
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/apps/
usr/share/icons/hicolor/128x128/apps/chromium.png
usr/share/icons/hicolor/16x16/
usr/share/icons/hicolor/16x16/apps/
usr/share/icons/hicolor/16x16/apps/chromium.png
usr/share/icons/hicolor/24x24/
usr/share/icons/hicolor/24x24/apps/
usr/share/icons/hicolor/24x24/apps/chromium.png
usr/share/icons/hicolor/256x256/
usr/share/icons/hicolor/256x256/apps/
usr/share/icons/hicolor/256x256/apps/chromium.png
usr/share/icons/hicolor/32x32/
usr/share/icons/hicolor/32x32/apps/
usr/share/icons/hicolor/32x32/apps/chromium.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/chromium.png
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/apps/
usr/share/icons/hicolor/64x64/apps/chromium.png
usr/share/licenses/
usr/share/licenses/chromium/
usr/share/licenses/chromium/LICENSE
usr/share/licenses/chromium/LICENSE.launcher
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/chromium.1.gz
usr/share/metainfo/
usr/share/metainfo/chromium.appdata.xml
