# files/qomui-git-0.8.2.r22.23650ab-2/files
%FILES%
usr/
usr/bin/
usr/bin/qomui-cli
usr/bin/qomui-gui
usr/bin/qomui-service
usr/lib/
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/qomui-0.8.2-py3.9.egg-info/
usr/lib/python3.9/site-packages/qomui-0.8.2-py3.9.egg-info/PKG-INFO
usr/lib/python3.9/site-packages/qomui-0.8.2-py3.9.egg-info/SOURCES.txt
usr/lib/python3.9/site-packages/qomui-0.8.2-py3.9.egg-info/dependency_links.txt
usr/lib/python3.9/site-packages/qomui-0.8.2-py3.9.egg-info/entry_points.txt
usr/lib/python3.9/site-packages/qomui-0.8.2-py3.9.egg-info/requires.txt
usr/lib/python3.9/site-packages/qomui-0.8.2-py3.9.egg-info/top_level.txt
usr/lib/python3.9/site-packages/qomui/
usr/lib/python3.9/site-packages/qomui/__init__.py
usr/lib/python3.9/site-packages/qomui/__pycache__/
usr/lib/python3.9/site-packages/qomui/__pycache__/__init__.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/bypass.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/bypass.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/config.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/config.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/dns_manager.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/dns_manager.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/firewall.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/firewall.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/latency.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/latency.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/monitor.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/monitor.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/mv_config.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/mv_config.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/plotter.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/plotter.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/profiles.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/profiles.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/qomui_cli.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/qomui_cli.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/qomui_gui.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/qomui_gui.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/qomui_service.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/qomui_service.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/tunnel.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/tunnel.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/update.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/update.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/utils.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/utils.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/widgets.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/qomui/__pycache__/widgets.cpython-39.pyc
usr/lib/python3.9/site-packages/qomui/bypass.py
usr/lib/python3.9/site-packages/qomui/config.py
usr/lib/python3.9/site-packages/qomui/dns_manager.py
usr/lib/python3.9/site-packages/qomui/firewall.py
usr/lib/python3.9/site-packages/qomui/latency.py
usr/lib/python3.9/site-packages/qomui/monitor.py
usr/lib/python3.9/site-packages/qomui/mv_config.py
usr/lib/python3.9/site-packages/qomui/plotter.py
usr/lib/python3.9/site-packages/qomui/profiles.py
usr/lib/python3.9/site-packages/qomui/qomui_cli.py
usr/lib/python3.9/site-packages/qomui/qomui_gui.py
usr/lib/python3.9/site-packages/qomui/qomui_service.py
usr/lib/python3.9/site-packages/qomui/tunnel.py
usr/lib/python3.9/site-packages/qomui/update.py
usr/lib/python3.9/site-packages/qomui/utils.py
usr/lib/python3.9/site-packages/qomui/widgets.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/qomui.service
usr/share/
usr/share/applications/
usr/share/applications/qomui.desktop
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.qomui.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.qomui.service.conf
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/qomui.svg
usr/share/icons/hicolor/scalable/apps/qomui_off.svg
usr/share/qomui/
usr/share/qomui/Airvpn_config
usr/share/qomui/Airvpn_config_old
usr/share/qomui/AzireVPN_config
usr/share/qomui/Mullvad_config
usr/share/qomui/Mullvad_config_old
usr/share/qomui/PIA_config
usr/share/qomui/PIA_config_old
usr/share/qomui/ProtonVPN_config
usr/share/qomui/ProtonVPN_config_old
usr/share/qomui/VERSION
usr/share/qomui/Windscribe_config
usr/share/qomui/Windscribe_config_old
usr/share/qomui/airvpn_api.pem
usr/share/qomui/airvpn_cacert.pem
usr/share/qomui/countries.json
usr/share/qomui/default_config.json
usr/share/qomui/firewall_default.json
usr/share/qomui/flags/
usr/share/qomui/flags/Afghanistan.png
usr/share/qomui/flags/Albania.png
usr/share/qomui/flags/Algeria.png
usr/share/qomui/flags/Andorra.png
usr/share/qomui/flags/Argentina.png
usr/share/qomui/flags/Armenia.png
usr/share/qomui/flags/Australia.png
usr/share/qomui/flags/Austria.png
usr/share/qomui/flags/Azerbaijan.png
usr/share/qomui/flags/Bahamas.png
usr/share/qomui/flags/Bangladesh.png
usr/share/qomui/flags/Belarus.png
usr/share/qomui/flags/Belgium.png
usr/share/qomui/flags/Bosnia and Herzegovina.png
usr/share/qomui/flags/Brazil.png
usr/share/qomui/flags/Bulgaria.png
usr/share/qomui/flags/Cameroon.png
usr/share/qomui/flags/Canada.png
usr/share/qomui/flags/Chile.png
usr/share/qomui/flags/China.png
usr/share/qomui/flags/Colombia.png
usr/share/qomui/flags/Costa Rica.png
usr/share/qomui/flags/Croatia.png
usr/share/qomui/flags/Cuba.png
usr/share/qomui/flags/Cyprus.png
usr/share/qomui/flags/Czech Republic.png
usr/share/qomui/flags/Denmark.png
usr/share/qomui/flags/Ecuador.png
usr/share/qomui/flags/Egypt.png
usr/share/qomui/flags/Estonia.png
usr/share/qomui/flags/Finland.png
usr/share/qomui/flags/France.png
usr/share/qomui/flags/Georgia.png
usr/share/qomui/flags/Germany.png
usr/share/qomui/flags/Greece.png
usr/share/qomui/flags/Hong Kong.png
usr/share/qomui/flags/Hungary.png
usr/share/qomui/flags/Iceland.png
usr/share/qomui/flags/India.png
usr/share/qomui/flags/Indonesia.png
usr/share/qomui/flags/Ireland.png
usr/share/qomui/flags/Israel.png
usr/share/qomui/flags/Italy.png
usr/share/qomui/flags/Jamaica.png
usr/share/qomui/flags/Japan.png
usr/share/qomui/flags/Latvia.png
usr/share/qomui/flags/Libya.png
usr/share/qomui/flags/Lithuania.png
usr/share/qomui/flags/Luxembourg.png
usr/share/qomui/flags/Macedonia, Republic of.png
usr/share/qomui/flags/Malaysia.png
usr/share/qomui/flags/Mexico.png
usr/share/qomui/flags/Moldova.png
usr/share/qomui/flags/Myanmar.png
usr/share/qomui/flags/Netherlands.png
usr/share/qomui/flags/New Zealand.png
usr/share/qomui/flags/Norway.png
usr/share/qomui/flags/Philippines.png
usr/share/qomui/flags/Poland.png
usr/share/qomui/flags/Portugal.png
usr/share/qomui/flags/Romania.png
usr/share/qomui/flags/Russia.png
usr/share/qomui/flags/Serbia.png
usr/share/qomui/flags/Singapore.png
usr/share/qomui/flags/Slovakia.png
usr/share/qomui/flags/Slovenia.png
usr/share/qomui/flags/South Africa.png
usr/share/qomui/flags/South Korea.png
usr/share/qomui/flags/Spain.png
usr/share/qomui/flags/Sweden.png
usr/share/qomui/flags/Switzerland.png
usr/share/qomui/flags/Taiwan.png
usr/share/qomui/flags/Thailand.png
usr/share/qomui/flags/Tunisia.png
usr/share/qomui/flags/Turkey.png
usr/share/qomui/flags/Ukraine.png
usr/share/qomui/flags/United Arab Emirates.png
usr/share/qomui/flags/United Kingdom.png
usr/share/qomui/flags/United States.png
usr/share/qomui/flags/Unknown.png
usr/share/qomui/flags/Viet Nam.png
usr/share/qomui/qomui.png
usr/share/qomui/scripts/
usr/share/qomui/scripts/bypass_route.sh
usr/share/qomui/scripts/bypass_up.sh
usr/share/qomui/scripts/hop.sh
usr/share/qomui/scripts/hop_down.sh
usr/share/qomui/ssl_config
