# files/ckb-next-0.4.4-5/files
%FILES%
usr/
usr/bin/
usr/bin/ckb-next
usr/bin/ckb-next-daemon
usr/bin/ckb-next-dev-detect
usr/bin/ckb-next.real
usr/lib/
usr/lib/ckb-next-animations/
usr/lib/ckb-next-animations/gradient
usr/lib/ckb-next-animations/heat
usr/lib/ckb-next-animations/invaders
usr/lib/ckb-next-animations/life
usr/lib/ckb-next-animations/pinwheel
usr/lib/ckb-next-animations/pipe
usr/lib/ckb-next-animations/rain
usr/lib/ckb-next-animations/random
usr/lib/ckb-next-animations/ripple
usr/lib/ckb-next-animations/snake
usr/lib/ckb-next-animations/wave
usr/lib/cmake/
usr/lib/cmake/ckb-next/
usr/lib/cmake/ckb-next/Animation/
usr/lib/cmake/ckb-next/Animation/AnimationConfig.cmake
usr/lib/cmake/ckb-next/Animation/AnimationConfigVersion.cmake
usr/lib/cmake/ckb-next/Animation/AnimationTargets.cmake
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ckb-next-daemon.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/99-ckb-next-daemon.rules
usr/share/
usr/share/applications/
usr/share/applications/ckb-next.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/apps/
usr/share/icons/hicolor/128x128/apps/ckb-next-monochrome.png
usr/share/icons/hicolor/128x128/apps/ckb-next.png
usr/share/icons/hicolor/16x16/
usr/share/icons/hicolor/16x16/apps/
usr/share/icons/hicolor/16x16/apps/ckb-next-monochrome.png
usr/share/icons/hicolor/16x16/apps/ckb-next.png
usr/share/icons/hicolor/22x22/
usr/share/icons/hicolor/22x22/apps/
usr/share/icons/hicolor/22x22/apps/ckb-next-monochrome.png
usr/share/icons/hicolor/22x22/apps/ckb-next.png
usr/share/icons/hicolor/256x256/
usr/share/icons/hicolor/256x256/apps/
usr/share/icons/hicolor/256x256/apps/ckb-next-monochrome.png
usr/share/icons/hicolor/256x256/apps/ckb-next.png
usr/share/icons/hicolor/32x32/
usr/share/icons/hicolor/32x32/apps/
usr/share/icons/hicolor/32x32/apps/ckb-next-monochrome.png
usr/share/icons/hicolor/32x32/apps/ckb-next.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/ckb-next-monochrome.png
usr/share/icons/hicolor/48x48/apps/ckb-next.png
usr/share/icons/hicolor/512x512/
usr/share/icons/hicolor/512x512/apps/
usr/share/icons/hicolor/512x512/apps/ckb-next-monochrome.png
usr/share/icons/hicolor/512x512/apps/ckb-next.png
