# files/tuxedo-control-center-1.0.14-1/files
%FILES%
etc/
etc/systemd/
etc/systemd/system/
etc/systemd/system/tccd-sleep.service
etc/systemd/system/tccd.service
opt/
opt/tuxedo-control-center/
opt/tuxedo-control-center/LICENSE.electron.txt
opt/tuxedo-control-center/LICENSES.chromium.html
opt/tuxedo-control-center/chrome-sandbox
opt/tuxedo-control-center/chrome_100_percent.pak
opt/tuxedo-control-center/chrome_200_percent.pak
opt/tuxedo-control-center/icudtl.dat
opt/tuxedo-control-center/libEGL.so
opt/tuxedo-control-center/libGLESv2.so
opt/tuxedo-control-center/libffmpeg.so
opt/tuxedo-control-center/libvk_swiftshader.so
opt/tuxedo-control-center/libvulkan.so
opt/tuxedo-control-center/locales/
opt/tuxedo-control-center/locales/am.pak
opt/tuxedo-control-center/locales/ar.pak
opt/tuxedo-control-center/locales/bg.pak
opt/tuxedo-control-center/locales/bn.pak
opt/tuxedo-control-center/locales/ca.pak
opt/tuxedo-control-center/locales/cs.pak
opt/tuxedo-control-center/locales/da.pak
opt/tuxedo-control-center/locales/de.pak
opt/tuxedo-control-center/locales/el.pak
opt/tuxedo-control-center/locales/en-GB.pak
opt/tuxedo-control-center/locales/en-US.pak
opt/tuxedo-control-center/locales/es-419.pak
opt/tuxedo-control-center/locales/es.pak
opt/tuxedo-control-center/locales/et.pak
opt/tuxedo-control-center/locales/fa.pak
opt/tuxedo-control-center/locales/fi.pak
opt/tuxedo-control-center/locales/fil.pak
opt/tuxedo-control-center/locales/fr.pak
opt/tuxedo-control-center/locales/gu.pak
opt/tuxedo-control-center/locales/he.pak
opt/tuxedo-control-center/locales/hi.pak
opt/tuxedo-control-center/locales/hr.pak
opt/tuxedo-control-center/locales/hu.pak
opt/tuxedo-control-center/locales/id.pak
opt/tuxedo-control-center/locales/it.pak
opt/tuxedo-control-center/locales/ja.pak
opt/tuxedo-control-center/locales/kn.pak
opt/tuxedo-control-center/locales/ko.pak
opt/tuxedo-control-center/locales/lt.pak
opt/tuxedo-control-center/locales/lv.pak
opt/tuxedo-control-center/locales/ml.pak
opt/tuxedo-control-center/locales/mr.pak
opt/tuxedo-control-center/locales/ms.pak
opt/tuxedo-control-center/locales/nb.pak
opt/tuxedo-control-center/locales/nl.pak
opt/tuxedo-control-center/locales/pl.pak
opt/tuxedo-control-center/locales/pt-BR.pak
opt/tuxedo-control-center/locales/pt-PT.pak
opt/tuxedo-control-center/locales/ro.pak
opt/tuxedo-control-center/locales/ru.pak
opt/tuxedo-control-center/locales/sk.pak
opt/tuxedo-control-center/locales/sl.pak
opt/tuxedo-control-center/locales/sr.pak
opt/tuxedo-control-center/locales/sv.pak
opt/tuxedo-control-center/locales/sw.pak
opt/tuxedo-control-center/locales/ta.pak
opt/tuxedo-control-center/locales/te.pak
opt/tuxedo-control-center/locales/th.pak
opt/tuxedo-control-center/locales/tr.pak
opt/tuxedo-control-center/locales/uk.pak
opt/tuxedo-control-center/locales/vi.pak
opt/tuxedo-control-center/locales/zh-CN.pak
opt/tuxedo-control-center/locales/zh-TW.pak
opt/tuxedo-control-center/resources.pak
opt/tuxedo-control-center/resources/
opt/tuxedo-control-center/resources/app.asar
opt/tuxedo-control-center/resources/dist/
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/dist-data/
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/dist-data/com.tuxedocomputers.tccd.conf
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/dist-data/de.tuxedocomputers.tcc.policy
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/dist-data/tccd-sleep.service
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/dist-data/tccd.service
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/dist-data/tuxedo-control-center-tray.desktop
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/dist-data/tuxedo-control-center.desktop
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/dist-data/tuxedo-control-center_256.svg
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/service/
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/service/TuxedoIOAPI.node
opt/tuxedo-control-center/resources/dist/tuxedo-control-center/data/service/tccd
opt/tuxedo-control-center/snapshot_blob.bin
opt/tuxedo-control-center/swiftshader/
opt/tuxedo-control-center/swiftshader/libEGL.so
opt/tuxedo-control-center/swiftshader/libGLESv2.so
opt/tuxedo-control-center/tuxedo-control-center
opt/tuxedo-control-center/v8_context_snapshot.bin
opt/tuxedo-control-center/vk_swiftshader_icd.json
usr/
usr/bin/
usr/bin/tuxedo-control-center
usr/share/
usr/share/applications/
usr/share/applications/tuxedo-control-center.desktop
usr/share/dbus-1/
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/com.tuxedocomputers.tccd.conf
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/apps/
usr/share/icons/hicolor/128x128/apps/tuxedo-control-center.png
usr/share/icons/hicolor/16x16/
usr/share/icons/hicolor/16x16/apps/
usr/share/icons/hicolor/16x16/apps/tuxedo-control-center.png
usr/share/icons/hicolor/24x24/
usr/share/icons/hicolor/24x24/apps/
usr/share/icons/hicolor/24x24/apps/tuxedo-control-center.png
usr/share/icons/hicolor/256x256/
usr/share/icons/hicolor/256x256/apps/
usr/share/icons/hicolor/256x256/apps/tuxedo-control-center.png
usr/share/icons/hicolor/32x32/
usr/share/icons/hicolor/32x32/apps/
usr/share/icons/hicolor/32x32/apps/tuxedo-control-center.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/tuxedo-control-center.png
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/apps/
usr/share/icons/hicolor/64x64/apps/tuxedo-control-center.png
usr/share/icons/hicolor/96x96/
usr/share/icons/hicolor/96x96/apps/
usr/share/icons/hicolor/96x96/apps/tuxedo-control-center.png
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/de.tuxedocomputers.tcc.policy
