# files/k3s-bin-1.24.4+k3s1-1/files
%FILES%
etc/
etc/systemd/
etc/systemd/system/
etc/systemd/system/k3s.service.env
usr/
usr/bin/
usr/bin/k3s
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/k3s.service
