# files/hsphfpd-git-r29.d294d06-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/org.hsphfpd.conf
usr/
usr/bin/
usr/bin/hsphfpd.pl
usr/bin/sco_features.pl
usr/bin/telephony_client.pl
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hsphfpd.service
usr/lib/systemd/user/
usr/lib/systemd/user/telephony_client.service
