# files/seafile-8.0.3-1/files
%FILES%
usr/
usr/bin/
usr/bin/seaf-cli
usr/bin/seaf-daemon
usr/include/
usr/include/seafile/
usr/include/seafile/seafile-error.h
usr/include/seafile/seafile-object.h
usr/include/seafile/seafile-rpc.h
usr/include/seafile/seafile.h
usr/lib/
usr/lib/libseafile.so
usr/lib/libseafile.so.0
usr/lib/libseafile.so.0.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/libseafile.pc
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/seafile/
usr/lib/python3.9/site-packages/seafile/__init__.py
usr/lib/python3.9/site-packages/seafile/__pycache__/
usr/lib/python3.9/site-packages/seafile/__pycache__/__init__.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/seafile/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/seafile/__pycache__/rpcclient.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/seafile/__pycache__/rpcclient.cpython-39.pyc
usr/lib/python3.9/site-packages/seafile/rpcclient.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/seaf-cli@.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/seaf-cli.1.gz
usr/share/man/man1/seaf-daemon.1.gz
