# files/openslp-2.0.0-2/files
%FILES%
etc/
etc/slp.conf
etc/slp.reg
etc/slp.spi
usr/
usr/bin/
usr/bin/slpd
usr/bin/slptool
usr/include/
usr/include/slp.h
usr/lib/
usr/lib/libslp.so
usr/lib/libslp.so.1
usr/lib/libslp.so.1.0.0
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/openslp.service
usr/share/
usr/share/doc/
usr/share/doc/openslp-2.0.0/
usr/share/doc/openslp-2.0.0/html/
usr/share/doc/openslp-2.0.0/html/IntroductionToSLP/
usr/share/doc/openslp-2.0.0/html/IntroductionToSLP/index.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/Callbacks.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/Divergence.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPAttrCallback.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPClose.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPDelAttrs.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPDereg.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPError.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPEscape.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPFindAttrs.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPFindScopes.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPFindSrvTypes.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPFindSrvs.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPFree.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPGetProperty.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPGetRefreshInterval.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPOpen.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPParseSrvURL.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPReg.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPRegReport.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPSetProperty.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPSrvTypeCallback.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPSrvURLCallback.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPTypes.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/SLPUnescape.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/Security.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/Syntax.html
usr/share/doc/openslp-2.0.0/html/ProgrammersGuide/index.html
usr/share/doc/openslp-2.0.0/html/UsersGuide/
usr/share/doc/openslp-2.0.0/html/UsersGuide/CommandLine.html
usr/share/doc/openslp-2.0.0/html/UsersGuide/FileLocations.html
usr/share/doc/openslp-2.0.0/html/UsersGuide/Installation.html
usr/share/doc/openslp-2.0.0/html/UsersGuide/Optimization.html
usr/share/doc/openslp-2.0.0/html/UsersGuide/Security.html
usr/share/doc/openslp-2.0.0/html/UsersGuide/SlpConf.html
usr/share/doc/openslp-2.0.0/html/UsersGuide/SlpReg.html
usr/share/doc/openslp-2.0.0/html/UsersGuide/WhenToRunSlpd.html
usr/share/doc/openslp-2.0.0/html/UsersGuide/WhoShouldRead.html
usr/share/doc/openslp-2.0.0/html/UsersGuide/index.html
usr/share/doc/openslp-2.0.0/html/faq.html
usr/share/doc/openslp-2.0.0/rfc/
usr/share/doc/openslp-2.0.0/rfc/rfc1766.txt
usr/share/doc/openslp-2.0.0/rfc/rfc2165.txt
usr/share/doc/openslp-2.0.0/rfc/rfc2254.txt
usr/share/doc/openslp-2.0.0/rfc/rfc2396.txt
usr/share/doc/openslp-2.0.0/rfc/rfc2608.txt
usr/share/doc/openslp-2.0.0/rfc/rfc2609.txt
usr/share/doc/openslp-2.0.0/rfc/rfc2610.txt
usr/share/doc/openslp-2.0.0/rfc/rfc2614.txt
usr/share/doc/openslp-2.0.0/rfc/rfc2926.txt
usr/share/doc/openslp-2.0.0/rfc/rfc3059.txt
usr/share/doc/openslp-2.0.0/rfc/rfc3082.txt
usr/share/doc/openslp-2.0.0/rfc/rfc3111.txt
usr/share/doc/openslp-2.0.0/rfc/rfc3224.txt
usr/share/doc/openslp-2.0.0/rfc/rfc3421.txt
usr/share/doc/openslp-2.0.0/rfc/rfc3528.txt
usr/share/doc/openslp-2.0.0/security/
usr/share/doc/openslp-2.0.0/security/openslp_security_whitepaper.html
usr/share/doc/openslp-2.0.0/security/srvreg-integrity.html
usr/share/doc/openslp-2.0.0/security/threat_analysis_min_security.html
usr/share/licenses/
usr/share/licenses/openslp/
usr/share/licenses/openslp/LICENSE
