# files/zramswap-5-1/files
%FILES%
etc/
etc/zramswap.conf
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/scripts/
usr/lib/systemd/scripts/zramctrl
usr/lib/systemd/system/
usr/lib/systemd/system/zramswap.service
