# files/shadowsocks-rust-git-r1610.b7387858-1/files
%FILES%
etc/
etc/shadowsocks/
etc/shadowsocks/config_ext_rust.json.example
etc/shadowsocks/config_rust.json.example
usr/
usr/bin/
usr/bin/sslocal-rust
usr/bin/ssserver-rust
usr/bin/ssurl-rust
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/shadowsocks-rust-server@.service
usr/lib/systemd/system/shadowsocks-rust@.service
