# files/syncplay-1.6.9-3/files
%FILES%
etc/
etc/syncplay/
etc/syncplay/server.conf
usr/
usr/bin/
usr/bin/syncplay
usr/bin/syncplay-server
usr/lib/
usr/lib/syncplay/
usr/lib/syncplay/syncplay/
usr/lib/syncplay/syncplay/__init__.py
usr/lib/syncplay/syncplay/client.py
usr/lib/syncplay/syncplay/clientManager.py
usr/lib/syncplay/syncplay/constants.py
usr/lib/syncplay/syncplay/ep_client.py
usr/lib/syncplay/syncplay/ep_server.py
usr/lib/syncplay/syncplay/messages.py
usr/lib/syncplay/syncplay/messages_de.py
usr/lib/syncplay/syncplay/messages_en.py
usr/lib/syncplay/syncplay/messages_es.py
usr/lib/syncplay/syncplay/messages_it.py
usr/lib/syncplay/syncplay/messages_pt_BR.py
usr/lib/syncplay/syncplay/messages_pt_PT.py
usr/lib/syncplay/syncplay/messages_ru.py
usr/lib/syncplay/syncplay/messages_tr.py
usr/lib/syncplay/syncplay/players/
usr/lib/syncplay/syncplay/players/__init__.py
usr/lib/syncplay/syncplay/players/basePlayer.py
usr/lib/syncplay/syncplay/players/iina.py
usr/lib/syncplay/syncplay/players/ipc_iina.py
usr/lib/syncplay/syncplay/players/mpc.py
usr/lib/syncplay/syncplay/players/mpcbe.py
usr/lib/syncplay/syncplay/players/mplayer.py
usr/lib/syncplay/syncplay/players/mpv.py
usr/lib/syncplay/syncplay/players/mpvnet.py
usr/lib/syncplay/syncplay/players/playerFactory.py
usr/lib/syncplay/syncplay/players/vlc.py
usr/lib/syncplay/syncplay/protocols.py
usr/lib/syncplay/syncplay/resources/
usr/lib/syncplay/syncplay/resources/IINA.png
usr/lib/syncplay/syncplay/resources/accept.png
usr/lib/syncplay/syncplay/resources/application_get.png
usr/lib/syncplay/syncplay/resources/arrow_refresh.png
usr/lib/syncplay/syncplay/resources/arrow_switch.png
usr/lib/syncplay/syncplay/resources/arrow_undo.png
usr/lib/syncplay/syncplay/resources/bullet_edit_centered.png
usr/lib/syncplay/syncplay/resources/bullet_right_grey.png
usr/lib/syncplay/syncplay/resources/chevrons_right.png
usr/lib/syncplay/syncplay/resources/clock_go.png
usr/lib/syncplay/syncplay/resources/cog.png
usr/lib/syncplay/syncplay/resources/cog_delete.png
usr/lib/syncplay/syncplay/resources/comments.png
usr/lib/syncplay/syncplay/resources/control_pause_blue.png
usr/lib/syncplay/syncplay/resources/control_play_blue.png
usr/lib/syncplay/syncplay/resources/cross.png
usr/lib/syncplay/syncplay/resources/cross_checkbox.png
usr/lib/syncplay/syncplay/resources/delete.png
usr/lib/syncplay/syncplay/resources/door_in.png
usr/lib/syncplay/syncplay/resources/door_open_edit.png
usr/lib/syncplay/syncplay/resources/email_go.png
usr/lib/syncplay/syncplay/resources/empty_checkbox.png
usr/lib/syncplay/syncplay/resources/error.png
usr/lib/syncplay/syncplay/resources/eye.png
usr/lib/syncplay/syncplay/resources/film_add.png
usr/lib/syncplay/syncplay/resources/film_edit.png
usr/lib/syncplay/syncplay/resources/film_folder_edit.png
usr/lib/syncplay/syncplay/resources/film_go.png
usr/lib/syncplay/syncplay/resources/film_link.png
usr/lib/syncplay/syncplay/resources/folder_explore.png
usr/lib/syncplay/syncplay/resources/folder_film.png
usr/lib/syncplay/syncplay/resources/help.png
usr/lib/syncplay/syncplay/resources/hicolor/
usr/lib/syncplay/syncplay/resources/hicolor/128x128/
usr/lib/syncplay/syncplay/resources/hicolor/128x128/apps/
usr/lib/syncplay/syncplay/resources/hicolor/128x128/apps/syncplay.png
usr/lib/syncplay/syncplay/resources/hicolor/16x16/
usr/lib/syncplay/syncplay/resources/hicolor/16x16/apps/
usr/lib/syncplay/syncplay/resources/hicolor/16x16/apps/syncplay.png
usr/lib/syncplay/syncplay/resources/hicolor/24x24/
usr/lib/syncplay/syncplay/resources/hicolor/24x24/apps/
usr/lib/syncplay/syncplay/resources/hicolor/24x24/apps/syncplay.png
usr/lib/syncplay/syncplay/resources/hicolor/256x256/
usr/lib/syncplay/syncplay/resources/hicolor/256x256/apps/
usr/lib/syncplay/syncplay/resources/hicolor/256x256/apps/syncplay.png
usr/lib/syncplay/syncplay/resources/hicolor/32x32/
usr/lib/syncplay/syncplay/resources/hicolor/32x32/apps/
usr/lib/syncplay/syncplay/resources/hicolor/32x32/apps/syncplay.png
usr/lib/syncplay/syncplay/resources/hicolor/48x48/
usr/lib/syncplay/syncplay/resources/hicolor/48x48/apps/
usr/lib/syncplay/syncplay/resources/hicolor/48x48/apps/syncplay.png
usr/lib/syncplay/syncplay/resources/hicolor/64x64/
usr/lib/syncplay/syncplay/resources/hicolor/64x64/apps/
usr/lib/syncplay/syncplay/resources/hicolor/64x64/apps/syncplay.png
usr/lib/syncplay/syncplay/resources/hicolor/96x96/
usr/lib/syncplay/syncplay/resources/hicolor/96x96/apps/
usr/lib/syncplay/syncplay/resources/hicolor/96x96/apps/syncplay.png
usr/lib/syncplay/syncplay/resources/house.png
usr/lib/syncplay/syncplay/resources/icon.icns
usr/lib/syncplay/syncplay/resources/icon.ico
usr/lib/syncplay/syncplay/resources/iina-bkg.png
usr/lib/syncplay/syncplay/resources/key_go.png
usr/lib/syncplay/syncplay/resources/license.rtf
usr/lib/syncplay/syncplay/resources/lock.png
usr/lib/syncplay/syncplay/resources/lock_green.png
usr/lib/syncplay/syncplay/resources/lock_green_dialog.png
usr/lib/syncplay/syncplay/resources/lock_green_dialog@2x.png
usr/lib/syncplay/syncplay/resources/lock_open.png
usr/lib/syncplay/syncplay/resources/lua/
usr/lib/syncplay/syncplay/resources/lua/intf/
usr/lib/syncplay/syncplay/resources/lua/intf/syncplay.lua
usr/lib/syncplay/syncplay/resources/macOS_dmg_bkg.tiff
usr/lib/syncplay/syncplay/resources/macOS_readme.pdf
usr/lib/syncplay/syncplay/resources/mpc-be.png
usr/lib/syncplay/syncplay/resources/mpc-hc.png
usr/lib/syncplay/syncplay/resources/mpc-hc64.png
usr/lib/syncplay/syncplay/resources/mplayer.png
usr/lib/syncplay/syncplay/resources/mpv.png
usr/lib/syncplay/syncplay/resources/mpvnet.png
usr/lib/syncplay/syncplay/resources/page_white_key.png
usr/lib/syncplay/syncplay/resources/shield_add.png
usr/lib/syncplay/syncplay/resources/shield_edit.png
usr/lib/syncplay/syncplay/resources/spinner.mng
usr/lib/syncplay/syncplay/resources/syncplay-server.desktop
usr/lib/syncplay/syncplay/resources/syncplay.desktop
usr/lib/syncplay/syncplay/resources/syncplay.png
usr/lib/syncplay/syncplay/resources/syncplayAbout.png
usr/lib/syncplay/syncplay/resources/syncplayAbout@2x.png
usr/lib/syncplay/syncplay/resources/syncplayintf.lua
usr/lib/syncplay/syncplay/resources/table_refresh.png
usr/lib/syncplay/syncplay/resources/third-party-notices.txt
usr/lib/syncplay/syncplay/resources/tick.png
usr/lib/syncplay/syncplay/resources/tick_checkbox.png
usr/lib/syncplay/syncplay/resources/timeline_marker.png
usr/lib/syncplay/syncplay/resources/user_comment.png
usr/lib/syncplay/syncplay/resources/user_key.png
usr/lib/syncplay/syncplay/resources/vlc.png
usr/lib/syncplay/syncplay/resources/world_add.png
usr/lib/syncplay/syncplay/resources/world_explore.png
usr/lib/syncplay/syncplay/resources/world_go.png
usr/lib/syncplay/syncplay/server.py
usr/lib/syncplay/syncplay/ui/
usr/lib/syncplay/syncplay/ui/ConfigurationGetter.py
usr/lib/syncplay/syncplay/ui/GuiConfiguration.py
usr/lib/syncplay/syncplay/ui/__init__.py
usr/lib/syncplay/syncplay/ui/consoleUI.py
usr/lib/syncplay/syncplay/ui/gui.py
usr/lib/syncplay/syncplay/utils.py
usr/lib/syncplay/syncplay/vendor/
usr/lib/syncplay/syncplay/vendor/Qt.py
usr/lib/syncplay/syncplay/vendor/__init__.py
usr/lib/syncplay/syncplay/vendor/darkdetect/
usr/lib/syncplay/syncplay/vendor/darkdetect/__init__.py
usr/lib/syncplay/syncplay/vendor/darkdetect/_detect.py
usr/lib/syncplay/syncplay/vendor/darkdetect/_dummy.py
usr/lib/syncplay/syncplay/vendor/python_mpv_jsonipc/
usr/lib/syncplay/syncplay/vendor/python_mpv_jsonipc/LICENSE.md
usr/lib/syncplay/syncplay/vendor/python_mpv_jsonipc/README.md
usr/lib/syncplay/syncplay/vendor/python_mpv_jsonipc/__init__.py
usr/lib/syncplay/syncplay/vendor/python_mpv_jsonipc/docs.md
usr/lib/syncplay/syncplay/vendor/python_mpv_jsonipc/python_mpv_jsonipc.py
usr/lib/syncplay/syncplay/vendor/python_mpv_jsonipc/setup.py
usr/lib/syncplay/syncplay/vendor/qt5reactor.py
usr/lib/syncplay/syncplayClient.py
usr/lib/syncplay/syncplayServer.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/syncplay@.service
usr/share/
usr/share/applications/
usr/share/applications/syncplay-server.desktop
usr/share/applications/syncplay.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/apps/
usr/share/icons/hicolor/128x128/apps/syncplay.png
usr/share/icons/hicolor/16x16/
usr/share/icons/hicolor/16x16/apps/
usr/share/icons/hicolor/16x16/apps/syncplay.png
usr/share/icons/hicolor/24x24/
usr/share/icons/hicolor/24x24/apps/
usr/share/icons/hicolor/24x24/apps/syncplay.png
usr/share/icons/hicolor/256x256/
usr/share/icons/hicolor/256x256/apps/
usr/share/icons/hicolor/256x256/apps/syncplay.png
usr/share/icons/hicolor/32x32/
usr/share/icons/hicolor/32x32/apps/
usr/share/icons/hicolor/32x32/apps/syncplay.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/syncplay.png
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/apps/
usr/share/icons/hicolor/64x64/apps/syncplay.png
usr/share/icons/hicolor/96x96/
usr/share/icons/hicolor/96x96/apps/
usr/share/icons/hicolor/96x96/apps/syncplay.png
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/syncplay-server.1.gz
usr/share/man/man1/syncplay.1.gz
usr/share/pixmaps/
usr/share/pixmaps/syncplay.png
