# files/zfs-utils-2.1.4-1/files
%FILES%
etc/
etc/default/
etc/default/zfs
etc/zfs/
etc/zfs/vdev_id.conf.alias.example
etc/zfs/vdev_id.conf.multipath.example
etc/zfs/vdev_id.conf.sas_direct.example
etc/zfs/vdev_id.conf.sas_switch.example
etc/zfs/vdev_id.conf.scsi.example
etc/zfs/zed.d/
etc/zfs/zed.d/all-syslog.sh
etc/zfs/zed.d/data-notify.sh
etc/zfs/zed.d/history_event-zfs-list-cacher.sh
etc/zfs/zed.d/pool_import-led.sh
etc/zfs/zed.d/resilver_finish-notify.sh
etc/zfs/zed.d/resilver_finish-start-scrub.sh
etc/zfs/zed.d/scrub_finish-notify.sh
etc/zfs/zed.d/statechange-led.sh
etc/zfs/zed.d/statechange-notify.sh
etc/zfs/zed.d/vdev_attach-led.sh
etc/zfs/zed.d/vdev_clear-led.sh
etc/zfs/zed.d/zed-functions.sh
etc/zfs/zed.d/zed.rc
etc/zfs/zfs-functions
etc/zfs/zpool.d/
etc/zfs/zpool.d/ata_err
etc/zfs/zpool.d/cmd_to
etc/zfs/zpool.d/defect
etc/zfs/zpool.d/dm-deps
etc/zfs/zpool.d/enc
etc/zfs/zpool.d/encdev
etc/zfs/zpool.d/fault_led
etc/zfs/zpool.d/health
etc/zfs/zpool.d/hours_on
etc/zfs/zpool.d/iostat
etc/zfs/zpool.d/iostat-10s
etc/zfs/zpool.d/iostat-1s
etc/zfs/zpool.d/label
etc/zfs/zpool.d/locate_led
etc/zfs/zpool.d/lsblk
etc/zfs/zpool.d/media
etc/zfs/zpool.d/model
etc/zfs/zpool.d/nonmed
etc/zfs/zpool.d/nvme_err
etc/zfs/zpool.d/off_ucor
etc/zfs/zpool.d/pend_sec
etc/zfs/zpool.d/pwr_cyc
etc/zfs/zpool.d/r_proc
etc/zfs/zpool.d/r_ucor
etc/zfs/zpool.d/realloc
etc/zfs/zpool.d/rep_ucor
etc/zfs/zpool.d/serial
etc/zfs/zpool.d/ses
etc/zfs/zpool.d/size
etc/zfs/zpool.d/slot
etc/zfs/zpool.d/smart
etc/zfs/zpool.d/smart_test
etc/zfs/zpool.d/smartx
etc/zfs/zpool.d/temp
etc/zfs/zpool.d/test_ended
etc/zfs/zpool.d/test_progress
etc/zfs/zpool.d/test_status
etc/zfs/zpool.d/test_type
etc/zfs/zpool.d/upath
etc/zfs/zpool.d/vendor
etc/zfs/zpool.d/w_proc
etc/zfs/zpool.d/w_ucor
usr/
usr/bin/
usr/bin/arc_summary
usr/bin/arcstat
usr/bin/dbufstat
usr/bin/fsck.zfs
usr/bin/mount.zfs
usr/bin/raidz_test
usr/bin/zdb
usr/bin/zed
usr/bin/zfs
usr/bin/zfs_ids_to_path
usr/bin/zgenhostid
usr/bin/zhack
usr/bin/zinject
usr/bin/zpool
usr/bin/zstream
usr/bin/zstreamdump
usr/bin/ztest
usr/bin/zvol_wait
usr/include/
usr/include/libspl/
usr/include/libspl/assert.h
usr/include/libspl/atomic.h
usr/include/libspl/ia32/
usr/include/libspl/ia32/sys/
usr/include/libspl/ia32/sys/asm_linkage.h
usr/include/libspl/libdevinfo.h
usr/include/libspl/libgen.h
usr/include/libspl/libshare.h
usr/include/libspl/limits.h
usr/include/libspl/locale.h
usr/include/libspl/rpc/
usr/include/libspl/rpc/xdr.h
usr/include/libspl/sscanf.h
usr/include/libspl/statcommon.h
usr/include/libspl/stdio.h
usr/include/libspl/stdlib.h
usr/include/libspl/string.h
usr/include/libspl/stropts.h
usr/include/libspl/sys/
usr/include/libspl/sys/acl.h
usr/include/libspl/sys/acl_impl.h
usr/include/libspl/sys/byteorder.h
usr/include/libspl/sys/callb.h
usr/include/libspl/sys/cmn_err.h
usr/include/libspl/sys/cred.h
usr/include/libspl/sys/debug.h
usr/include/libspl/sys/dkio.h
usr/include/libspl/sys/dklabel.h
usr/include/libspl/sys/dktp/
usr/include/libspl/sys/dktp/fdisk.h
usr/include/libspl/sys/errno.h
usr/include/libspl/sys/feature_tests.h
usr/include/libspl/sys/int_limits.h
usr/include/libspl/sys/int_types.h
usr/include/libspl/sys/inttypes.h
usr/include/libspl/sys/isa_defs.h
usr/include/libspl/sys/kmem.h
usr/include/libspl/sys/kstat.h
usr/include/libspl/sys/list.h
usr/include/libspl/sys/list_impl.h
usr/include/libspl/sys/mhd.h
usr/include/libspl/sys/mkdev.h
usr/include/libspl/sys/mnttab.h
usr/include/libspl/sys/mount.h
usr/include/libspl/sys/param.h
usr/include/libspl/sys/policy.h
usr/include/libspl/sys/poll.h
usr/include/libspl/sys/priv.h
usr/include/libspl/sys/processor.h
usr/include/libspl/sys/sha2.h
usr/include/libspl/sys/simd.h
usr/include/libspl/sys/stack.h
usr/include/libspl/sys/stat.h
usr/include/libspl/sys/stdtypes.h
usr/include/libspl/sys/strings.h
usr/include/libspl/sys/stropts.h
usr/include/libspl/sys/sunddi.h
usr/include/libspl/sys/sysmacros.h
usr/include/libspl/sys/systeminfo.h
usr/include/libspl/sys/time.h
usr/include/libspl/sys/trace_spl.h
usr/include/libspl/sys/trace_zfs.h
usr/include/libspl/sys/types.h
usr/include/libspl/sys/types32.h
usr/include/libspl/sys/tzfile.h
usr/include/libspl/sys/uio.h
usr/include/libspl/sys/va_list.h
usr/include/libspl/sys/varargs.h
usr/include/libspl/sys/vnode.h
usr/include/libspl/sys/vtoc.h
usr/include/libspl/sys/wmsum.h
usr/include/libspl/sys/zfs_context_os.h
usr/include/libspl/sys/zone.h
usr/include/libspl/thread.h
usr/include/libspl/tzfile.h
usr/include/libspl/ucred.h
usr/include/libspl/umem.h
usr/include/libspl/unistd.h
usr/include/libspl/zone.h
usr/include/libzfs/
usr/include/libzfs/cityhash.h
usr/include/libzfs/libnvpair.h
usr/include/libzfs/libuutil.h
usr/include/libzfs/libuutil_common.h
usr/include/libzfs/libuutil_impl.h
usr/include/libzfs/libzfs.h
usr/include/libzfs/libzfs_core.h
usr/include/libzfs/libzfs_impl.h
usr/include/libzfs/libzfsbootenv.h
usr/include/libzfs/libzutil.h
usr/include/libzfs/sys/
usr/include/libzfs/sys/abd.h
usr/include/libzfs/sys/abd_impl.h
usr/include/libzfs/sys/aggsum.h
usr/include/libzfs/sys/arc.h
usr/include/libzfs/sys/arc_impl.h
usr/include/libzfs/sys/avl.h
usr/include/libzfs/sys/avl_impl.h
usr/include/libzfs/sys/bitops.h
usr/include/libzfs/sys/blkptr.h
usr/include/libzfs/sys/bplist.h
usr/include/libzfs/sys/bpobj.h
usr/include/libzfs/sys/bptree.h
usr/include/libzfs/sys/bqueue.h
usr/include/libzfs/sys/btree.h
usr/include/libzfs/sys/crypto/
usr/include/libzfs/sys/crypto/api.h
usr/include/libzfs/sys/crypto/common.h
usr/include/libzfs/sys/crypto/icp.h
usr/include/libzfs/sys/dataset_kstats.h
usr/include/libzfs/sys/dbuf.h
usr/include/libzfs/sys/ddt.h
usr/include/libzfs/sys/dmu.h
usr/include/libzfs/sys/dmu_impl.h
usr/include/libzfs/sys/dmu_objset.h
usr/include/libzfs/sys/dmu_recv.h
usr/include/libzfs/sys/dmu_redact.h
usr/include/libzfs/sys/dmu_send.h
usr/include/libzfs/sys/dmu_traverse.h
usr/include/libzfs/sys/dmu_tx.h
usr/include/libzfs/sys/dmu_zfetch.h
usr/include/libzfs/sys/dnode.h
usr/include/libzfs/sys/dsl_bookmark.h
usr/include/libzfs/sys/dsl_crypt.h
usr/include/libzfs/sys/dsl_dataset.h
usr/include/libzfs/sys/dsl_deadlist.h
usr/include/libzfs/sys/dsl_deleg.h
usr/include/libzfs/sys/dsl_destroy.h
usr/include/libzfs/sys/dsl_dir.h
usr/include/libzfs/sys/dsl_pool.h
usr/include/libzfs/sys/dsl_prop.h
usr/include/libzfs/sys/dsl_scan.h
usr/include/libzfs/sys/dsl_synctask.h
usr/include/libzfs/sys/dsl_userhold.h
usr/include/libzfs/sys/edonr.h
usr/include/libzfs/sys/efi_partition.h
usr/include/libzfs/sys/fm/
usr/include/libzfs/sys/fm/fs/
usr/include/libzfs/sys/fm/fs/zfs.h
usr/include/libzfs/sys/fm/protocol.h
usr/include/libzfs/sys/fm/util.h
usr/include/libzfs/sys/frame.h
usr/include/libzfs/sys/fs/
usr/include/libzfs/sys/fs/zfs.h
usr/include/libzfs/sys/hkdf.h
usr/include/libzfs/sys/lua/
usr/include/libzfs/sys/lua/lauxlib.h
usr/include/libzfs/sys/lua/lua.h
usr/include/libzfs/sys/lua/luaconf.h
usr/include/libzfs/sys/lua/lualib.h
usr/include/libzfs/sys/metaslab.h
usr/include/libzfs/sys/metaslab_impl.h
usr/include/libzfs/sys/mmp.h
usr/include/libzfs/sys/mntent.h
usr/include/libzfs/sys/mod.h
usr/include/libzfs/sys/multilist.h
usr/include/libzfs/sys/note.h
usr/include/libzfs/sys/nvpair.h
usr/include/libzfs/sys/nvpair_impl.h
usr/include/libzfs/sys/objlist.h
usr/include/libzfs/sys/pathname.h
usr/include/libzfs/sys/qat.h
usr/include/libzfs/sys/range_tree.h
usr/include/libzfs/sys/rrwlock.h
usr/include/libzfs/sys/sa.h
usr/include/libzfs/sys/sa_impl.h
usr/include/libzfs/sys/skein.h
usr/include/libzfs/sys/spa.h
usr/include/libzfs/sys/spa_boot.h
usr/include/libzfs/sys/spa_checkpoint.h
usr/include/libzfs/sys/spa_checksum.h
usr/include/libzfs/sys/spa_impl.h
usr/include/libzfs/sys/spa_log_spacemap.h
usr/include/libzfs/sys/space_map.h
usr/include/libzfs/sys/space_reftree.h
usr/include/libzfs/sys/sysevent.h
usr/include/libzfs/sys/sysevent/
usr/include/libzfs/sys/sysevent/dev.h
usr/include/libzfs/sys/sysevent/eventdefs.h
usr/include/libzfs/sys/txg.h
usr/include/libzfs/sys/txg_impl.h
usr/include/libzfs/sys/u8_textprep.h
usr/include/libzfs/sys/u8_textprep_data.h
usr/include/libzfs/sys/uberblock.h
usr/include/libzfs/sys/uberblock_impl.h
usr/include/libzfs/sys/uio_impl.h
usr/include/libzfs/sys/unique.h
usr/include/libzfs/sys/uuid.h
usr/include/libzfs/sys/vdev.h
usr/include/libzfs/sys/vdev_disk.h
usr/include/libzfs/sys/vdev_draid.h
usr/include/libzfs/sys/vdev_file.h
usr/include/libzfs/sys/vdev_impl.h
usr/include/libzfs/sys/vdev_indirect_births.h
usr/include/libzfs/sys/vdev_indirect_mapping.h
usr/include/libzfs/sys/vdev_initialize.h
usr/include/libzfs/sys/vdev_raidz.h
usr/include/libzfs/sys/vdev_raidz_impl.h
usr/include/libzfs/sys/vdev_rebuild.h
usr/include/libzfs/sys/vdev_removal.h
usr/include/libzfs/sys/vdev_trim.h
usr/include/libzfs/sys/xvattr.h
usr/include/libzfs/sys/zap.h
usr/include/libzfs/sys/zap_impl.h
usr/include/libzfs/sys/zap_leaf.h
usr/include/libzfs/sys/zcp.h
usr/include/libzfs/sys/zcp_global.h
usr/include/libzfs/sys/zcp_iter.h
usr/include/libzfs/sys/zcp_prop.h
usr/include/libzfs/sys/zcp_set.h
usr/include/libzfs/sys/zfeature.h
usr/include/libzfs/sys/zfs_acl.h
usr/include/libzfs/sys/zfs_bootenv.h
usr/include/libzfs/sys/zfs_context.h
usr/include/libzfs/sys/zfs_debug.h
usr/include/libzfs/sys/zfs_delay.h
usr/include/libzfs/sys/zfs_file.h
usr/include/libzfs/sys/zfs_fuid.h
usr/include/libzfs/sys/zfs_project.h
usr/include/libzfs/sys/zfs_quota.h
usr/include/libzfs/sys/zfs_racct.h
usr/include/libzfs/sys/zfs_ratelimit.h
usr/include/libzfs/sys/zfs_refcount.h
usr/include/libzfs/sys/zfs_rlock.h
usr/include/libzfs/sys/zfs_sa.h
usr/include/libzfs/sys/zfs_stat.h
usr/include/libzfs/sys/zfs_sysfs.h
usr/include/libzfs/sys/zfs_vfsops.h
usr/include/libzfs/sys/zfs_vnops.h
usr/include/libzfs/sys/zfs_znode.h
usr/include/libzfs/sys/zil.h
usr/include/libzfs/sys/zil_impl.h
usr/include/libzfs/sys/zio.h
usr/include/libzfs/sys/zio_checksum.h
usr/include/libzfs/sys/zio_compress.h
usr/include/libzfs/sys/zio_crypt.h
usr/include/libzfs/sys/zio_impl.h
usr/include/libzfs/sys/zio_priority.h
usr/include/libzfs/sys/zrlock.h
usr/include/libzfs/sys/zstd/
usr/include/libzfs/sys/zstd/zstd.h
usr/include/libzfs/sys/zthr.h
usr/include/libzfs/thread_pool.h
usr/include/libzfs/zfeature_common.h
usr/include/libzfs/zfs_comutil.h
usr/include/libzfs/zfs_deleg.h
usr/include/libzfs/zfs_fletcher.h
usr/include/libzfs/zfs_namecheck.h
usr/include/libzfs/zfs_prop.h
usr/lib/
usr/lib/dracut/
usr/lib/dracut/modules.d/
usr/lib/dracut/modules.d/02zfsexpandknowledge/
usr/lib/dracut/modules.d/02zfsexpandknowledge/module-setup.sh
usr/lib/dracut/modules.d/90zfs/
usr/lib/dracut/modules.d/90zfs/export-zfs.sh
usr/lib/dracut/modules.d/90zfs/import-opts-generator.sh
usr/lib/dracut/modules.d/90zfs/module-setup.sh
usr/lib/dracut/modules.d/90zfs/mount-zfs.sh
usr/lib/dracut/modules.d/90zfs/parse-zfs.sh
usr/lib/dracut/modules.d/90zfs/zfs-env-bootfs.service
usr/lib/dracut/modules.d/90zfs/zfs-generator.sh
usr/lib/dracut/modules.d/90zfs/zfs-lib.sh
usr/lib/dracut/modules.d/90zfs/zfs-load-key.sh
usr/lib/dracut/modules.d/90zfs/zfs-needshutdown.sh
usr/lib/dracut/modules.d/90zfs/zfs-rollback-bootfs.service
usr/lib/dracut/modules.d/90zfs/zfs-snapshot-bootfs.service
usr/lib/initcpio/
usr/lib/initcpio/hooks/
usr/lib/initcpio/hooks/zfs
usr/lib/initcpio/install/
usr/lib/initcpio/install/zfs
usr/lib/libnvpair.so
usr/lib/libnvpair.so.3
usr/lib/libnvpair.so.3.0.0
usr/lib/libuutil.so
usr/lib/libuutil.so.3
usr/lib/libuutil.so.3.0.0
usr/lib/libzfs.so
usr/lib/libzfs.so.4
usr/lib/libzfs.so.4.1.0
usr/lib/libzfs_core.so
usr/lib/libzfs_core.so.3
usr/lib/libzfs_core.so.3.0.0
usr/lib/libzfsbootenv.so
usr/lib/libzfsbootenv.so.1
usr/lib/libzfsbootenv.so.1.0.0
usr/lib/libzpool.so
usr/lib/libzpool.so.5
usr/lib/libzpool.so.5.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/libzfs.pc
usr/lib/pkgconfig/libzfs_core.pc
usr/lib/pkgconfig/libzfsbootenv.pc
usr/lib/security/
usr/lib/security/pam_zfs_key.so
usr/lib/systemd/
usr/lib/systemd/system-generators/
usr/lib/systemd/system-generators/zfs-mount-generator
usr/lib/systemd/system-preset/
usr/lib/systemd/system-preset/50-zfs.preset
usr/lib/systemd/system/
usr/lib/systemd/system/zfs-import-cache.service
usr/lib/systemd/system/zfs-import-scan.service
usr/lib/systemd/system/zfs-import.service
usr/lib/systemd/system/zfs-import.target
usr/lib/systemd/system/zfs-mount.service
usr/lib/systemd/system/zfs-scrub-monthly@.timer
usr/lib/systemd/system/zfs-scrub-weekly@.timer
usr/lib/systemd/system/zfs-scrub@.service
usr/lib/systemd/system/zfs-share.service
usr/lib/systemd/system/zfs-volume-wait.service
usr/lib/systemd/system/zfs-volumes.target
usr/lib/systemd/system/zfs-zed.service
usr/lib/systemd/system/zfs.target
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-zvol.rules
usr/lib/udev/rules.d/69-vdev.rules
usr/lib/udev/rules.d/90-zfs.rules
usr/lib/udev/vdev_id
usr/lib/udev/zvol_id
usr/lib/zfs/
usr/lib/zfs/zfs/
usr/lib/zfs/zfs/zed.d/
usr/lib/zfs/zfs/zed.d/all-debug.sh
usr/lib/zfs/zfs/zed.d/all-syslog.sh
usr/lib/zfs/zfs/zed.d/data-notify.sh
usr/lib/zfs/zfs/zed.d/generic-notify.sh
usr/lib/zfs/zfs/zed.d/history_event-zfs-list-cacher.sh
usr/lib/zfs/zfs/zed.d/pool_import-led.sh
usr/lib/zfs/zfs/zed.d/resilver_finish-notify.sh
usr/lib/zfs/zfs/zed.d/resilver_finish-start-scrub.sh
usr/lib/zfs/zfs/zed.d/scrub_finish-notify.sh
usr/lib/zfs/zfs/zed.d/statechange-led.sh
usr/lib/zfs/zfs/zed.d/statechange-notify.sh
usr/lib/zfs/zfs/zed.d/trim_finish-notify.sh
usr/lib/zfs/zfs/zed.d/vdev_attach-led.sh
usr/lib/zfs/zfs/zed.d/vdev_clear-led.sh
usr/lib/zfs/zfs/zpool.d/
usr/lib/zfs/zfs/zpool.d/ata_err
usr/lib/zfs/zfs/zpool.d/cmd_to
usr/lib/zfs/zfs/zpool.d/defect
usr/lib/zfs/zfs/zpool.d/dm-deps
usr/lib/zfs/zfs/zpool.d/enc
usr/lib/zfs/zfs/zpool.d/encdev
usr/lib/zfs/zfs/zpool.d/fault_led
usr/lib/zfs/zfs/zpool.d/health
usr/lib/zfs/zfs/zpool.d/hours_on
usr/lib/zfs/zfs/zpool.d/iostat
usr/lib/zfs/zfs/zpool.d/iostat-10s
usr/lib/zfs/zfs/zpool.d/iostat-1s
usr/lib/zfs/zfs/zpool.d/label
usr/lib/zfs/zfs/zpool.d/locate_led
usr/lib/zfs/zfs/zpool.d/lsblk
usr/lib/zfs/zfs/zpool.d/media
usr/lib/zfs/zfs/zpool.d/model
usr/lib/zfs/zfs/zpool.d/nonmed
usr/lib/zfs/zfs/zpool.d/nvme_err
usr/lib/zfs/zfs/zpool.d/off_ucor
usr/lib/zfs/zfs/zpool.d/pend_sec
usr/lib/zfs/zfs/zpool.d/pwr_cyc
usr/lib/zfs/zfs/zpool.d/r_proc
usr/lib/zfs/zfs/zpool.d/r_ucor
usr/lib/zfs/zfs/zpool.d/realloc
usr/lib/zfs/zfs/zpool.d/rep_ucor
usr/lib/zfs/zfs/zpool.d/serial
usr/lib/zfs/zfs/zpool.d/ses
usr/lib/zfs/zfs/zpool.d/size
usr/lib/zfs/zfs/zpool.d/slot
usr/lib/zfs/zfs/zpool.d/smart
usr/lib/zfs/zfs/zpool.d/smart_test
usr/lib/zfs/zfs/zpool.d/smartx
usr/lib/zfs/zfs/zpool.d/temp
usr/lib/zfs/zfs/zpool.d/test_ended
usr/lib/zfs/zfs/zpool.d/test_progress
usr/lib/zfs/zfs/zpool.d/test_status
usr/lib/zfs/zfs/zpool.d/test_type
usr/lib/zfs/zfs/zpool.d/upath
usr/lib/zfs/zfs/zpool.d/vendor
usr/lib/zfs/zfs/zpool.d/w_proc
usr/lib/zfs/zfs/zpool.d/w_ucor
usr/lib/zfs/zfs/zpool_influxdb
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/zfs
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/arcstat.1.gz
usr/share/man/man1/raidz_test.1.gz
usr/share/man/man1/test-runner.1.gz
usr/share/man/man1/zhack.1.gz
usr/share/man/man1/ztest.1.gz
usr/share/man/man1/zvol_wait.1.gz
usr/share/man/man4/
usr/share/man/man4/spl.4.gz
usr/share/man/man4/zfs.4.gz
usr/share/man/man5/
usr/share/man/man5/vdev_id.conf.5.gz
usr/share/man/man7/
usr/share/man/man7/zfsconcepts.7.gz
usr/share/man/man7/zfsprops.7.gz
usr/share/man/man7/zpool-features.7.gz
usr/share/man/man7/zpoolconcepts.7.gz
usr/share/man/man7/zpoolprops.7.gz
usr/share/man/man8/
usr/share/man/man8/fsck.zfs.8.gz
usr/share/man/man8/mount.zfs.8.gz
usr/share/man/man8/vdev_id.8.gz
usr/share/man/man8/zdb.8.gz
usr/share/man/man8/zed.8.gz
usr/share/man/man8/zfs-allow.8.gz
usr/share/man/man8/zfs-bookmark.8.gz
usr/share/man/man8/zfs-change-key.8.gz
usr/share/man/man8/zfs-clone.8.gz
usr/share/man/man8/zfs-create.8.gz
usr/share/man/man8/zfs-destroy.8.gz
usr/share/man/man8/zfs-diff.8.gz
usr/share/man/man8/zfs-get.8.gz
usr/share/man/man8/zfs-groupspace.8.gz
usr/share/man/man8/zfs-hold.8.gz
usr/share/man/man8/zfs-inherit.8.gz
usr/share/man/man8/zfs-jail.8.gz
usr/share/man/man8/zfs-list.8.gz
usr/share/man/man8/zfs-load-key.8.gz
usr/share/man/man8/zfs-mount-generator.8.gz
usr/share/man/man8/zfs-mount.8.gz
usr/share/man/man8/zfs-program.8.gz
usr/share/man/man8/zfs-project.8.gz
usr/share/man/man8/zfs-projectspace.8.gz
usr/share/man/man8/zfs-promote.8.gz
usr/share/man/man8/zfs-receive.8.gz
usr/share/man/man8/zfs-recv.8.gz
usr/share/man/man8/zfs-redact.8.gz
usr/share/man/man8/zfs-release.8.gz
usr/share/man/man8/zfs-rename.8.gz
usr/share/man/man8/zfs-rollback.8.gz
usr/share/man/man8/zfs-send.8.gz
usr/share/man/man8/zfs-set.8.gz
usr/share/man/man8/zfs-share.8.gz
usr/share/man/man8/zfs-snapshot.8.gz
usr/share/man/man8/zfs-unallow.8.gz
usr/share/man/man8/zfs-unjail.8.gz
usr/share/man/man8/zfs-unload-key.8.gz
usr/share/man/man8/zfs-unmount.8.gz
usr/share/man/man8/zfs-upgrade.8.gz
usr/share/man/man8/zfs-userspace.8.gz
usr/share/man/man8/zfs-wait.8.gz
usr/share/man/man8/zfs.8.gz
usr/share/man/man8/zfs_ids_to_path.8.gz
usr/share/man/man8/zgenhostid.8.gz
usr/share/man/man8/zinject.8.gz
usr/share/man/man8/zpool-add.8.gz
usr/share/man/man8/zpool-attach.8.gz
usr/share/man/man8/zpool-checkpoint.8.gz
usr/share/man/man8/zpool-clear.8.gz
usr/share/man/man8/zpool-create.8.gz
usr/share/man/man8/zpool-destroy.8.gz
usr/share/man/man8/zpool-detach.8.gz
usr/share/man/man8/zpool-events.8.gz
usr/share/man/man8/zpool-export.8.gz
usr/share/man/man8/zpool-get.8.gz
usr/share/man/man8/zpool-history.8.gz
usr/share/man/man8/zpool-import.8.gz
usr/share/man/man8/zpool-initialize.8.gz
usr/share/man/man8/zpool-iostat.8.gz
usr/share/man/man8/zpool-labelclear.8.gz
usr/share/man/man8/zpool-list.8.gz
usr/share/man/man8/zpool-offline.8.gz
usr/share/man/man8/zpool-online.8.gz
usr/share/man/man8/zpool-reguid.8.gz
usr/share/man/man8/zpool-remove.8.gz
usr/share/man/man8/zpool-reopen.8.gz
usr/share/man/man8/zpool-replace.8.gz
usr/share/man/man8/zpool-resilver.8.gz
usr/share/man/man8/zpool-scrub.8.gz
usr/share/man/man8/zpool-set.8.gz
usr/share/man/man8/zpool-split.8.gz
usr/share/man/man8/zpool-status.8.gz
usr/share/man/man8/zpool-sync.8.gz
usr/share/man/man8/zpool-trim.8.gz
usr/share/man/man8/zpool-upgrade.8.gz
usr/share/man/man8/zpool-wait.8.gz
usr/share/man/man8/zpool.8.gz
usr/share/man/man8/zpool_influxdb.8.gz
usr/share/man/man8/zstream.8.gz
usr/share/man/man8/zstreamdump.8.gz
usr/share/pam-configs/
usr/share/pam-configs/zfs_key
