# files/windscribe-cli-1.4_51-1/files
%FILES%
etc/
etc/bash_completion.d/
etc/bash_completion.d/windscribe_complete
etc/windscribe/
usr/
usr/bin/
usr/bin/windscribe
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/windscribe.service
usr/share/
usr/share/licenses/
usr/share/licenses/windscribe-cli/
usr/share/licenses/windscribe-cli/LICENSE
