# files/tblock-2.0.0-1/files
%FILES%
etc/
etc/tblock.conf
usr/
usr/bin/
usr/bin/tblock
usr/bin/tblockc
usr/bin/tblockd
usr/lib/
usr/lib/python3.10/
usr/lib/python3.10/site-packages/
usr/lib/python3.10/site-packages/tblock-2.0.0-py3.10.egg-info/
usr/lib/python3.10/site-packages/tblock-2.0.0-py3.10.egg-info/PKG-INFO
usr/lib/python3.10/site-packages/tblock-2.0.0-py3.10.egg-info/SOURCES.txt
usr/lib/python3.10/site-packages/tblock-2.0.0-py3.10.egg-info/dependency_links.txt
usr/lib/python3.10/site-packages/tblock-2.0.0-py3.10.egg-info/entry_points.txt
usr/lib/python3.10/site-packages/tblock-2.0.0-py3.10.egg-info/requires.txt
usr/lib/python3.10/site-packages/tblock-2.0.0-py3.10.egg-info/top_level.txt
usr/lib/python3.10/site-packages/tblock/
usr/lib/python3.10/site-packages/tblock/__init__.py
usr/lib/python3.10/site-packages/tblock/__main__.py
usr/lib/python3.10/site-packages/tblock/__pycache__/
usr/lib/python3.10/site-packages/tblock/__pycache__/__init__.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/__main__.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/__main__.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/cli.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/cli.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/compat.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/compat.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/config.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/config.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/const.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/const.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/daemon.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/daemon.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/exceptions.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/exceptions.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/filters.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/filters.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/hosts.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/hosts.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/rules.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/rules.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/style.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/style.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/utils.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/__pycache__/utils.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/argumentor/
usr/lib/python3.10/site-packages/tblock/argumentor/__init__.py
usr/lib/python3.10/site-packages/tblock/argumentor/__pycache__/
usr/lib/python3.10/site-packages/tblock/argumentor/__pycache__/__init__.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/argumentor/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/argumentor/__pycache__/exceptions.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/argumentor/__pycache__/exceptions.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/argumentor/exceptions.py
usr/lib/python3.10/site-packages/tblock/cli.py
usr/lib/python3.10/site-packages/tblock/compat.py
usr/lib/python3.10/site-packages/tblock/config.py
usr/lib/python3.10/site-packages/tblock/const.py
usr/lib/python3.10/site-packages/tblock/converter/
usr/lib/python3.10/site-packages/tblock/converter/__init__.py
usr/lib/python3.10/site-packages/tblock/converter/__main__.py
usr/lib/python3.10/site-packages/tblock/converter/__pycache__/
usr/lib/python3.10/site-packages/tblock/converter/__pycache__/__init__.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/converter/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/converter/__pycache__/__main__.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/converter/__pycache__/__main__.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/converter/__pycache__/syntaxtools.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/tblock/converter/__pycache__/syntaxtools.cpython-310.pyc
usr/lib/python3.10/site-packages/tblock/converter/syntaxtools.py
usr/lib/python3.10/site-packages/tblock/daemon.py
usr/lib/python3.10/site-packages/tblock/exceptions.py
usr/lib/python3.10/site-packages/tblock/filters.py
usr/lib/python3.10/site-packages/tblock/hosts.py
usr/lib/python3.10/site-packages/tblock/rules.py
usr/lib/python3.10/site-packages/tblock/style.py
usr/lib/python3.10/site-packages/tblock/utils.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tblockd.service
usr/share/
usr/share/doc/
usr/share/doc/tblock/
usr/share/doc/tblock/README.txt
usr/share/licenses/
usr/share/licenses/tblock/
usr/share/licenses/tblock/LICENSE.txt
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/tblock.1.gz
usr/share/man/man1/tblockc.1.gz
usr/share/man/man1/tblockd.1.gz
usr/share/man/man5/
usr/share/man/man5/tblock.conf.5.gz
