# files/expressvpn-3.13.0.8_1-1/files
%FILES%
usr/
usr/bin/
usr/bin/expressvpn
usr/bin/expressvpn-agent
usr/bin/expressvpn-browser-helper
usr/bin/expressvpnd
usr/lib/
usr/lib/expressvpn/
usr/lib/expressvpn/chrome/
usr/lib/expressvpn/chrome/com.expressvpn.helper.json
usr/lib/expressvpn/expressvpn-agent.desktop
usr/lib/expressvpn/firefox/
usr/lib/expressvpn/firefox/com.expressvpn.helper.json
usr/lib/expressvpn/icon.png
usr/lib/expressvpn/libxvclient.so
usr/lib/expressvpn/lightway
usr/lib/expressvpn/openvpn
usr/lib/expressvpn/version
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/expressvpn.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/expressvpn
usr/share/doc/
usr/share/doc/expressvpn/
usr/share/doc/expressvpn/COPYRIGHT
usr/share/doc/expressvpn/changelog.gz
usr/share/licenses/
usr/share/licenses/expressvpn/
usr/share/licenses/expressvpn/COPYRIGHT
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/expressvpn.1.gz
var/
var/lib/
var/lib/expressvpn/
var/lib/expressvpn/certs/
