# files/asusctl-git-4.0.7.r16.gbac2ba6-1/files
%FILES%
etc/
etc/asusd/
etc/asusd/asusd-ledmodes.toml
usr/
usr/bin/
usr/bin/asus-notify
usr/bin/asusctl
usr/bin/asusd
usr/bin/asusd-user
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/asusd.service
usr/lib/systemd/user/
usr/lib/systemd/user/asus-notify.service
usr/lib/systemd/user/asusd-user.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/99-asusd.rules
usr/share/
usr/share/asusd/
usr/share/asusd/anime/
usr/share/asusd/anime/asus/
usr/share/asusd/anime/asus/festive/
usr/share/asusd/anime/asus/festive/Cupid.gif
usr/share/asusd/anime/asus/festive/Firework.gif
usr/share/asusd/anime/asus/festive/Halloween.gif
usr/share/asusd/anime/asus/festive/Happy Holiday.gif
usr/share/asusd/anime/asus/festive/Happy new year.gif
usr/share/asusd/anime/asus/festive/Lantern.gif
usr/share/asusd/anime/asus/festive/Love u mom.gif
usr/share/asusd/anime/asus/festive/Mother's day.gif
usr/share/asusd/anime/asus/festive/Valentine's Day.gif
usr/share/asusd/anime/asus/festive/Year of the Ox.gif
usr/share/asusd/anime/asus/gaming/
usr/share/asusd/anime/asus/gaming/Bird.gif
usr/share/asusd/anime/asus/gaming/Controller.gif
usr/share/asusd/anime/asus/gaming/FPS.gif
usr/share/asusd/anime/asus/gaming/Fight.gif
usr/share/asusd/anime/asus/gaming/Keyboard.gif
usr/share/asusd/anime/asus/gaming/MOBA.gif
usr/share/asusd/anime/asus/gaming/UFO.gif
usr/share/asusd/anime/asus/music/
usr/share/asusd/anime/asus/music/DJ.gif
usr/share/asusd/anime/asus/music/Diamond.gif
usr/share/asusd/anime/asus/music/Music-player.gif
usr/share/asusd/anime/asus/rog/
usr/share/asusd/anime/asus/rog/For-those-who-dare.gif
usr/share/asusd/anime/asus/rog/For-those-who-dare_2.gif
usr/share/asusd/anime/asus/rog/Fragment.gif
usr/share/asusd/anime/asus/rog/Infinite-triangle.gif
usr/share/asusd/anime/asus/rog/Kaleidoscope1.gif
usr/share/asusd/anime/asus/rog/Kaleidoscope2.gif
usr/share/asusd/anime/asus/rog/Kaleidoscope2.png
usr/share/asusd/anime/asus/rog/ROG city.gif
usr/share/asusd/anime/asus/rog/ROG glitch.gif
usr/share/asusd/anime/asus/rog/Sunset.gif
usr/share/asusd/anime/asus/trend/
usr/share/asusd/anime/asus/trend/Dog.gif
usr/share/asusd/anime/asus/trend/Hero.gif
usr/share/asusd/anime/asus/trend/Ski.gif
usr/share/asusd/anime/asus/trend/The scream.gif
usr/share/asusd/anime/asus/trend/Wave.gif
usr/share/asusd/anime/custom/
usr/share/asusd/anime/custom/nyancat_zombie.gif
usr/share/asusd/anime/custom/rust.png
usr/share/asusd/anime/custom/sonic-run.gif
usr/share/asusd/anime/custom/sonic-wait.gif
usr/share/dbus-1/
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/asusd.conf
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/asusctl.fish
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/512x512/
usr/share/icons/hicolor/512x512/apps/
usr/share/icons/hicolor/512x512/apps/asus_notif_green.png
usr/share/icons/hicolor/512x512/apps/asus_notif_red.png
usr/share/icons/hicolor/512x512/apps/asus_notif_yellow.png
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/status/
usr/share/icons/hicolor/scalable/status/gpu-compute.svg
usr/share/icons/hicolor/scalable/status/gpu-hybrid.svg
usr/share/icons/hicolor/scalable/status/gpu-integrated.svg
usr/share/icons/hicolor/scalable/status/gpu-nvidia.svg
usr/share/icons/hicolor/scalable/status/gpu-vfio.svg
usr/share/icons/hicolor/scalable/status/notification-reboot.svg
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_asusctl
