# files/tgt-1.0.80-1/files
%FILES%
etc/
etc/tgt/
etc/tgt/conf.d/
etc/tgt/examples/
etc/tgt/examples/targets.conf.example
etc/tgt/examples/targets.conf.vtl.L700
etc/tgt/examples/targets.conf.vtl.MSL2024
etc/tgt/targets.conf
usr/
usr/bin/
usr/bin/tgt-admin
usr/bin/tgt-setup-lun
usr/bin/tgtadm
usr/bin/tgtd
usr/bin/tgtimg
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tgtd.service
usr/share/
usr/share/doc/
usr/share/doc/tgt/
usr/share/doc/tgt/README.iscsi
usr/share/doc/tgt/README.iser
usr/share/doc/tgt/README.lu_configuration
usr/share/doc/tgt/README.mmc
usr/share/doc/tgt/README.rbd
usr/share/doc/tgt/html/
usr/share/doc/tgt/html/targets.conf.5.html
usr/share/doc/tgt/html/tgt-admin.8.html
usr/share/doc/tgt/html/tgt-setup-lun.8.html
usr/share/doc/tgt/html/tgtadm.8.html
usr/share/doc/tgt/html/tgtd.8.html
usr/share/doc/tgt/html/tgtimg.8.html
usr/share/doc/tgt/tmf.txt
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/targets.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/tgt-admin.8.gz
usr/share/man/man8/tgt-setup-lun.8.gz
usr/share/man/man8/tgtadm.8.gz
usr/share/man/man8/tgtd.8.gz
usr/share/man/man8/tgtimg.8.gz
