# files/nvidia-470xx-utils-470.129.06-1/files
%FILES%
usr/
usr/bin/
usr/bin/nvidia-bug-report.sh
usr/bin/nvidia-cuda-mps-control
usr/bin/nvidia-cuda-mps-server
usr/bin/nvidia-debugdump
usr/bin/nvidia-modprobe
usr/bin/nvidia-ngx-updater
usr/bin/nvidia-persistenced
usr/bin/nvidia-sleep.sh
usr/bin/nvidia-smi
usr/bin/nvidia-xconfig
usr/lib/
usr/lib/firmware/
usr/lib/firmware/nvidia/
usr/lib/firmware/nvidia/470.129.06/
usr/lib/firmware/nvidia/470.129.06/gsp.bin
usr/lib/libEGL_nvidia.so
usr/lib/libEGL_nvidia.so.0
usr/lib/libEGL_nvidia.so.470.129.06
usr/lib/libGLESv1_CM_nvidia.so
usr/lib/libGLESv1_CM_nvidia.so.1
usr/lib/libGLESv1_CM_nvidia.so.470.129.06
usr/lib/libGLESv2_nvidia.so
usr/lib/libGLESv2_nvidia.so.2
usr/lib/libGLESv2_nvidia.so.470.129.06
usr/lib/libGLX_nvidia.so
usr/lib/libGLX_nvidia.so.0
usr/lib/libGLX_nvidia.so.470.129.06
usr/lib/libcuda.so
usr/lib/libcuda.so.1
usr/lib/libcuda.so.470.129.06
usr/lib/libnvcuvid.so
usr/lib/libnvcuvid.so.1
usr/lib/libnvcuvid.so.470.129.06
usr/lib/libnvidia-allocator.so
usr/lib/libnvidia-allocator.so.1
usr/lib/libnvidia-allocator.so.470.129.06
usr/lib/libnvidia-cbl.so
usr/lib/libnvidia-cbl.so.470.129.06
usr/lib/libnvidia-cfg.so
usr/lib/libnvidia-cfg.so.1
usr/lib/libnvidia-cfg.so.470.129.06
usr/lib/libnvidia-eglcore.so
usr/lib/libnvidia-eglcore.so.470.129.06
usr/lib/libnvidia-encode.so
usr/lib/libnvidia-encode.so.1
usr/lib/libnvidia-encode.so.470.129.06
usr/lib/libnvidia-fbc.so
usr/lib/libnvidia-fbc.so.1
usr/lib/libnvidia-fbc.so.470.129.06
usr/lib/libnvidia-glcore.so
usr/lib/libnvidia-glcore.so.470.129.06
usr/lib/libnvidia-glsi.so
usr/lib/libnvidia-glsi.so.470.129.06
usr/lib/libnvidia-glvkspirv.so
usr/lib/libnvidia-glvkspirv.so.470.129.06
usr/lib/libnvidia-ifr.so
usr/lib/libnvidia-ifr.so.1
usr/lib/libnvidia-ifr.so.470.129.06
usr/lib/libnvidia-ml.so
usr/lib/libnvidia-ml.so.1
usr/lib/libnvidia-ml.so.470.129.06
usr/lib/libnvidia-ngx.so
usr/lib/libnvidia-ngx.so.1
usr/lib/libnvidia-ngx.so.470.129.06
usr/lib/libnvidia-opticalflow.so
usr/lib/libnvidia-opticalflow.so.1
usr/lib/libnvidia-opticalflow.so.470.129.06
usr/lib/libnvidia-ptxjitcompiler.so
usr/lib/libnvidia-ptxjitcompiler.so.1
usr/lib/libnvidia-ptxjitcompiler.so.470.129.06
usr/lib/libnvidia-rtcore.so
usr/lib/libnvidia-rtcore.so.470.129.06
usr/lib/libnvidia-tls.so
usr/lib/libnvidia-tls.so.470.129.06
usr/lib/libnvidia-vulkan-producer.so
usr/lib/libnvidia-vulkan-producer.so.1
usr/lib/libnvidia-vulkan-producer.so.470.129.06
usr/lib/libnvoptix.so
usr/lib/libnvoptix.so.1
usr/lib/libnvoptix.so.470.129.06
usr/lib/modprobe.d/
usr/lib/modprobe.d/nvidia-470xx-utils.conf
usr/lib/modules-load.d/
usr/lib/modules-load.d/nvidia-470xx-utils.conf
usr/lib/nvidia/
usr/lib/nvidia/wine/
usr/lib/nvidia/wine/_nvngx.dll
usr/lib/nvidia/wine/nvngx.dll
usr/lib/nvidia/xorg/
usr/lib/nvidia/xorg/libglxserver_nvidia.so
usr/lib/nvidia/xorg/libglxserver_nvidia.so.1
usr/lib/nvidia/xorg/libglxserver_nvidia.so.470.129.06
usr/lib/systemd/
usr/lib/systemd/system-sleep/
usr/lib/systemd/system-sleep/nvidia
usr/lib/systemd/system/
usr/lib/systemd/system/nvidia-hibernate.service
usr/lib/systemd/system/nvidia-persistenced.service
usr/lib/systemd/system/nvidia-resume.service
usr/lib/systemd/system/nvidia-suspend.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/nvidia-470xx-utils.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-nvidia-470xx.rules
usr/lib/vdpau/
usr/lib/vdpau/libvdpau_nvidia.so
usr/lib/vdpau/libvdpau_nvidia.so.1
usr/lib/vdpau/libvdpau_nvidia.so.470.129.06
usr/lib/xorg/
usr/lib/xorg/modules/
usr/lib/xorg/modules/drivers/
usr/lib/xorg/modules/drivers/nvidia_drv.so
usr/share/
usr/share/X11/
usr/share/X11/xorg.conf.d/
usr/share/X11/xorg.conf.d/10-nvidia-drm-outputclass.conf
usr/share/doc/
usr/share/doc/nvidia-utils
usr/share/doc/nvidia/
usr/share/doc/nvidia/NVIDIA_Changelog
usr/share/doc/nvidia/README
usr/share/doc/nvidia/html/
usr/share/doc/nvidia/html/acknowledgements.html
usr/share/doc/nvidia/html/addressingcapabilities.html
usr/share/doc/nvidia/html/addtlresources.html
usr/share/doc/nvidia/html/appendices.html
usr/share/doc/nvidia/html/audiosupport.html
usr/share/doc/nvidia/html/commonproblems.html
usr/share/doc/nvidia/html/configlaptop.html
usr/share/doc/nvidia/html/configmultxscreens.html
usr/share/doc/nvidia/html/configtwinview.html
usr/share/doc/nvidia/html/depth30.html
usr/share/doc/nvidia/html/displaydevicenames.html
usr/share/doc/nvidia/html/dma_issues.html
usr/share/doc/nvidia/html/dpi.html
usr/share/doc/nvidia/html/dynamicpowermanagement.html
usr/share/doc/nvidia/html/editxconfig.html
usr/share/doc/nvidia/html/egpu.html
usr/share/doc/nvidia/html/faq.html
usr/share/doc/nvidia/html/flippingubb.html
usr/share/doc/nvidia/html/framelock.html
usr/share/doc/nvidia/html/glxsupport.html
usr/share/doc/nvidia/html/gpunames.html
usr/share/doc/nvidia/html/gsp.html
usr/share/doc/nvidia/html/i2c.html
usr/share/doc/nvidia/html/index.html
usr/share/doc/nvidia/html/installationandconfiguration.html
usr/share/doc/nvidia/html/installdriver.html
usr/share/doc/nvidia/html/installedcomponents.html
usr/share/doc/nvidia/html/introduction.html
usr/share/doc/nvidia/html/kms.html
usr/share/doc/nvidia/html/knownissues.html
usr/share/doc/nvidia/html/minimumrequirements.html
usr/share/doc/nvidia/html/newusertips.html
usr/share/doc/nvidia/html/ngx.html
usr/share/doc/nvidia/html/nvidia-debugdump.html
usr/share/doc/nvidia/html/nvidia-ml.html
usr/share/doc/nvidia/html/nvidia-peermem.html
usr/share/doc/nvidia/html/nvidia-persistenced.html
usr/share/doc/nvidia/html/nvidia-smi.html
usr/share/doc/nvidia/html/nvidiasettings.html
usr/share/doc/nvidia/html/openglenvvariables.html
usr/share/doc/nvidia/html/optimus.html
usr/share/doc/nvidia/html/powermanagement.html
usr/share/doc/nvidia/html/primerenderoffload.html
usr/share/doc/nvidia/html/procinterface.html
usr/share/doc/nvidia/html/profiles.html
usr/share/doc/nvidia/html/programmingmodes.html
usr/share/doc/nvidia/html/randr14.html
usr/share/doc/nvidia/html/retpoline.html
usr/share/doc/nvidia/html/selectdriver.html
usr/share/doc/nvidia/html/sli.html
usr/share/doc/nvidia/html/supportedchips.html
usr/share/doc/nvidia/html/vdpausupport.html
usr/share/doc/nvidia/html/xcompositeextension.html
usr/share/doc/nvidia/html/xconfigoptions.html
usr/share/doc/nvidia/html/xineramaglx.html
usr/share/doc/nvidia/html/xrandrextension.html
usr/share/doc/nvidia/html/xwayland.html
usr/share/glvnd/
usr/share/glvnd/egl_vendor.d/
usr/share/glvnd/egl_vendor.d/10_nvidia.json
usr/share/licenses/
usr/share/licenses/nvidia-utils/
usr/share/licenses/nvidia-utils/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/nvidia-cuda-mps-control.1.gz
usr/share/man/man1/nvidia-modprobe.1.gz
usr/share/man/man1/nvidia-persistenced.1.gz
usr/share/man/man1/nvidia-smi.1.gz
usr/share/man/man1/nvidia-xconfig.1.gz
usr/share/nvidia/
usr/share/nvidia/nvidia-application-profiles-470.129.06-key-documentation
usr/share/nvidia/nvidia-application-profiles-470.129.06-rc
usr/share/vulkan/
usr/share/vulkan/icd.d/
usr/share/vulkan/icd.d/nvidia_icd.json
usr/share/vulkan/implicit_layer.d/
usr/share/vulkan/implicit_layer.d/nvidia_layers.json
