# files/garuda-live-systemd-r24.ecd887a-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/garuda-live-keyring.service
usr/lib/systemd/system/garuda-live.service
usr/lib/systemd/system/mhwd-live-net.service
usr/lib/systemd/system/mhwd-live.service
