# files/garuda-system-maintenance-git-r16.bb6b34a-1/files
%FILES%
etc/
etc/systemd/
etc/systemd/system/
etc/systemd/system/garuda-system-maintenance@.service
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/garuda-system-maintenance-git.desktop
usr/
usr/bin/
usr/bin/garuda-system-maintenance
usr/lib/
usr/lib/garuda-system-maintenance/
usr/lib/garuda-system-maintenance/update-packages
usr/share/
usr/share/applications/
usr/share/applications/garuda-system-maintenance-settings.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/garuda-system-maintenance.svg
usr/share/knotifications5/
usr/share/knotifications5/garuda-system-maintenance.notifyrc
usr/share/polkit-1/
usr/share/polkit-1/rules.d/
usr/share/polkit-1/rules.d/garuda-system-maintenance.rules
