# files/pigpio-79-1/files
%FILES%
usr/
usr/bin/
usr/bin/pig2vcd
usr/bin/pigpiod
usr/bin/pigs
usr/include/
usr/include/pigpio.h
usr/include/pigpiod_if.h
usr/include/pigpiod_if2.h
usr/lib/
usr/lib/libpigpio.so
usr/lib/libpigpio.so.1
usr/lib/libpigpiod_if.so
usr/lib/libpigpiod_if.so.1
usr/lib/libpigpiod_if2.so
usr/lib/libpigpiod_if2.so.1
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/__pycache__/
usr/lib/python3.9/site-packages/__pycache__/pigpio.cpython-39.pyc
usr/lib/python3.9/site-packages/pigpio-1.78-py3.9.egg-info
usr/lib/python3.9/site-packages/pigpio.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pigpiod.service
usr/share/
usr/share/licenses/
usr/share/licenses/pigpio/
usr/share/licenses/pigpio/UNLICENCE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/pig2vcd.1.gz
usr/share/man/man1/pigpiod.1.gz
usr/share/man/man1/pigs.1.gz
usr/share/man/man3/
usr/share/man/man3/pigpio.3.gz
usr/share/man/man3/pigpiod_if.3.gz
usr/share/man/man3/pigpiod_if2.3.gz
