# files/tpacpi-bat-git-3.1.r12.g4959b52-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/tpacpi
usr/
usr/bin/
usr/bin/tpacpi-bat
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tpacpi-bat.service
