# files/libsdrplay-3.07.1-1/files
%FILES%
etc/
etc/systemd/
etc/systemd/system/
etc/systemd/system/sdrplay.service
etc/udev/
etc/udev/rules.d/
etc/udev/rules.d/66-mirics.rules
usr/
usr/bin/
usr/bin/sdrplay_apiService
usr/include/
usr/include/sdrplay_api.h
usr/include/sdrplay_api_callback.h
usr/include/sdrplay_api_control.h
usr/include/sdrplay_api_dev.h
usr/include/sdrplay_api_rsp1a.h
usr/include/sdrplay_api_rsp2.h
usr/include/sdrplay_api_rspDuo.h
usr/include/sdrplay_api_rspDx.h
usr/include/sdrplay_api_rx_channel.h
usr/include/sdrplay_api_tuner.h
usr/lib/
usr/lib/libsdrplay_api.so
usr/lib/libsdrplay_api.so.2
usr/lib/libsdrplay_api.so.3.07
usr/share/
usr/share/licenses/
usr/share/licenses/libsdrplay/
usr/share/licenses/libsdrplay/LICENSE
