# files/android-sdk-platform-tools-33.0.2-1/files
%FILES%
etc/
etc/profile.d/
etc/profile.d/android-sdk-platform-tools.csh
etc/profile.d/android-sdk-platform-tools.sh
opt/
opt/android-sdk/
opt/android-sdk/platform-tools/
opt/android-sdk/platform-tools/NOTICE.txt
opt/android-sdk/platform-tools/adb
opt/android-sdk/platform-tools/dmtracedump
opt/android-sdk/platform-tools/e2fsdroid
opt/android-sdk/platform-tools/etc1tool
opt/android-sdk/platform-tools/fastboot
opt/android-sdk/platform-tools/hprof-conv
opt/android-sdk/platform-tools/lib64/
opt/android-sdk/platform-tools/lib64/libc++.so
opt/android-sdk/platform-tools/make_f2fs
opt/android-sdk/platform-tools/make_f2fs_casefold
opt/android-sdk/platform-tools/mke2fs
opt/android-sdk/platform-tools/mke2fs.conf
opt/android-sdk/platform-tools/sload_f2fs
opt/android-sdk/platform-tools/source.properties
opt/android-sdk/platform-tools/sqlite3
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/adb.service
usr/share/
usr/share/licenses/
usr/share/licenses/android-sdk-platform-tools/
usr/share/licenses/android-sdk-platform-tools/license.html
