# files/open-fprintd-0.6-3/files
%FILES%
usr/
usr/lib/
usr/lib/open-fprintd/
usr/lib/open-fprintd/open-fprintd
usr/lib/open-fprintd/resume.py
usr/lib/open-fprintd/suspend.py
usr/lib/python3.10/
usr/lib/python3.10/site-packages/
usr/lib/python3.10/site-packages/open_fprintd-0.6-py3.10.egg-info/
usr/lib/python3.10/site-packages/open_fprintd-0.6-py3.10.egg-info/PKG-INFO
usr/lib/python3.10/site-packages/open_fprintd-0.6-py3.10.egg-info/SOURCES.txt
usr/lib/python3.10/site-packages/open_fprintd-0.6-py3.10.egg-info/dependency_links.txt
usr/lib/python3.10/site-packages/open_fprintd-0.6-py3.10.egg-info/top_level.txt
usr/lib/python3.10/site-packages/openfprintd/
usr/lib/python3.10/site-packages/openfprintd/__init__.py
usr/lib/python3.10/site-packages/openfprintd/__pycache__/
usr/lib/python3.10/site-packages/openfprintd/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/openfprintd/__pycache__/device.cpython-310.pyc
usr/lib/python3.10/site-packages/openfprintd/__pycache__/manager.cpython-310.pyc
usr/lib/python3.10/site-packages/openfprintd/device.py
usr/lib/python3.10/site-packages/openfprintd/manager.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/open-fprintd-resume.service
usr/lib/systemd/system/open-fprintd-suspend.service
usr/lib/systemd/system/open-fprintd.service
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/net.reactivated.Fprint.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/net.reactivated.Fprint.conf
