# files/mullvad-vpn-beta-2022.1.stable1-1/files
%FILES%
opt/
opt/Mullvad VPN/
opt/Mullvad VPN/LICENSE.electron.txt
opt/Mullvad VPN/LICENSES.chromium.html
opt/Mullvad VPN/chrome-sandbox
opt/Mullvad VPN/chrome_100_percent.pak
opt/Mullvad VPN/chrome_200_percent.pak
opt/Mullvad VPN/chrome_crashpad_handler
opt/Mullvad VPN/icudtl.dat
opt/Mullvad VPN/libEGL.so
opt/Mullvad VPN/libGLESv2.so
opt/Mullvad VPN/libffmpeg.so
opt/Mullvad VPN/libvk_swiftshader.so
opt/Mullvad VPN/libvulkan.so.1
opt/Mullvad VPN/locales/
opt/Mullvad VPN/locales/am.pak
opt/Mullvad VPN/locales/ar.pak
opt/Mullvad VPN/locales/bg.pak
opt/Mullvad VPN/locales/bn.pak
opt/Mullvad VPN/locales/ca.pak
opt/Mullvad VPN/locales/cs.pak
opt/Mullvad VPN/locales/da.pak
opt/Mullvad VPN/locales/de.pak
opt/Mullvad VPN/locales/el.pak
opt/Mullvad VPN/locales/en-GB.pak
opt/Mullvad VPN/locales/en-US.pak
opt/Mullvad VPN/locales/es-419.pak
opt/Mullvad VPN/locales/es.pak
opt/Mullvad VPN/locales/et.pak
opt/Mullvad VPN/locales/fa.pak
opt/Mullvad VPN/locales/fi.pak
opt/Mullvad VPN/locales/fil.pak
opt/Mullvad VPN/locales/fr.pak
opt/Mullvad VPN/locales/gu.pak
opt/Mullvad VPN/locales/he.pak
opt/Mullvad VPN/locales/hi.pak
opt/Mullvad VPN/locales/hr.pak
opt/Mullvad VPN/locales/hu.pak
opt/Mullvad VPN/locales/id.pak
opt/Mullvad VPN/locales/it.pak
opt/Mullvad VPN/locales/ja.pak
opt/Mullvad VPN/locales/kn.pak
opt/Mullvad VPN/locales/ko.pak
opt/Mullvad VPN/locales/lt.pak
opt/Mullvad VPN/locales/lv.pak
opt/Mullvad VPN/locales/ml.pak
opt/Mullvad VPN/locales/mr.pak
opt/Mullvad VPN/locales/ms.pak
opt/Mullvad VPN/locales/nb.pak
opt/Mullvad VPN/locales/nl.pak
opt/Mullvad VPN/locales/pl.pak
opt/Mullvad VPN/locales/pt-BR.pak
opt/Mullvad VPN/locales/pt-PT.pak
opt/Mullvad VPN/locales/ro.pak
opt/Mullvad VPN/locales/ru.pak
opt/Mullvad VPN/locales/sk.pak
opt/Mullvad VPN/locales/sl.pak
opt/Mullvad VPN/locales/sr.pak
opt/Mullvad VPN/locales/sv.pak
opt/Mullvad VPN/locales/sw.pak
opt/Mullvad VPN/locales/ta.pak
opt/Mullvad VPN/locales/te.pak
opt/Mullvad VPN/locales/th.pak
opt/Mullvad VPN/locales/tr.pak
opt/Mullvad VPN/locales/uk.pak
opt/Mullvad VPN/locales/vi.pak
opt/Mullvad VPN/locales/zh-CN.pak
opt/Mullvad VPN/locales/zh-TW.pak
opt/Mullvad VPN/mullvad-gui
opt/Mullvad VPN/mullvad-vpn
opt/Mullvad VPN/resources.pak
opt/Mullvad VPN/resources/
opt/Mullvad VPN/resources/CHANGELOG.md
opt/Mullvad VPN/resources/api-ip-address.txt
opt/Mullvad VPN/resources/app.asar
opt/Mullvad VPN/resources/ca.crt
opt/Mullvad VPN/resources/libtalpid_openvpn_plugin.so
opt/Mullvad VPN/resources/mullvad-daemon
opt/Mullvad VPN/resources/mullvad-daemon.conf
opt/Mullvad VPN/resources/mullvad-daemon.service
opt/Mullvad VPN/resources/mullvad-problem-report
opt/Mullvad VPN/resources/mullvad-setup
opt/Mullvad VPN/resources/openvpn
opt/Mullvad VPN/resources/relays.json
opt/Mullvad VPN/resources/sslocal
opt/Mullvad VPN/snapshot_blob.bin
opt/Mullvad VPN/swiftshader/
opt/Mullvad VPN/swiftshader/libEGL.so
opt/Mullvad VPN/swiftshader/libGLESv2.so
opt/Mullvad VPN/v8_context_snapshot.bin
opt/Mullvad VPN/vk_swiftshader_icd.json
usr/
usr/bin/
usr/bin/mullvad
usr/bin/mullvad-exclude
usr/bin/mullvad-problem-report
usr/bin/mullvad-vpn
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mullvad-daemon.service
usr/share/
usr/share/applications/
usr/share/applications/mullvad-vpn.desktop
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/mullvad
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/mullvad.fish
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/1024x1024/
usr/share/icons/hicolor/1024x1024/apps/
usr/share/icons/hicolor/1024x1024/apps/mullvad-vpn.png
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/apps/
usr/share/icons/hicolor/128x128/apps/mullvad-vpn.png
usr/share/icons/hicolor/16x16/
usr/share/icons/hicolor/16x16/apps/
usr/share/icons/hicolor/16x16/apps/mullvad-vpn.png
usr/share/icons/hicolor/256x256/
usr/share/icons/hicolor/256x256/apps/
usr/share/icons/hicolor/256x256/apps/mullvad-vpn.png
usr/share/icons/hicolor/32x32/
usr/share/icons/hicolor/32x32/apps/
usr/share/icons/hicolor/32x32/apps/mullvad-vpn.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/mullvad-vpn.png
usr/share/icons/hicolor/512x512/
usr/share/icons/hicolor/512x512/apps/
usr/share/icons/hicolor/512x512/apps/mullvad-vpn.png
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/apps/
usr/share/icons/hicolor/64x64/apps/mullvad-vpn.png
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_mullvad
