# files/headscale-0.15.0-1/files
%FILES%
etc/
etc/headscale/
etc/headscale/config.yaml
usr/
usr/bin/
usr/bin/headscale
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/headscale.service
usr/share/
usr/share/licenses/
usr/share/licenses/headscale/
usr/share/licenses/headscale/LICENSE
