# files/openrgb-git-r2409.bcc89af4-1/files
%FILES%
usr/
usr/bin/
usr/bin/openrgb
usr/lib/
usr/lib/modules-load.d/
usr/lib/modules-load.d/openrgb.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/openrgb.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-openrgb.rules
usr/share/
usr/share/applications/
usr/share/applications/OpenRGB.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/apps/
usr/share/icons/hicolor/128x128/apps/OpenRGB.png
usr/share/metainfo/
usr/share/metainfo/org.openrgb.OpenRGB.metainfo.xml
