# files/das-keyboard-q-3.3.3-1/files
%FILES%
etc/
etc/systemd/
etc/systemd/system/
etc/systemd/system/das_keyboard_q-service.service
etc/udev/
etc/udev/rules.d/
etc/udev/rules.d/99-daskeyboard.rules
usr/
usr/bin/
usr/bin/das-keyboard-q
usr/bin/das_keyboard_q-service
usr/lib/
usr/lib/das-keyboard-q/
usr/lib/das-keyboard-q/LICENSES.chromium.html
usr/lib/das-keyboard-q/chrome-sandbox
usr/lib/das-keyboard-q/chrome_100_percent.pak
usr/lib/das-keyboard-q/chrome_200_percent.pak
usr/lib/das-keyboard-q/das-keyboard-q
usr/lib/das-keyboard-q/icudtl.dat
usr/lib/das-keyboard-q/libEGL.so
usr/lib/das-keyboard-q/libGLESv2.so
usr/lib/das-keyboard-q/libffmpeg.so
usr/lib/das-keyboard-q/libvk_swiftshader.so
usr/lib/das-keyboard-q/libvulkan.so.1
usr/lib/das-keyboard-q/locales/
usr/lib/das-keyboard-q/locales/am.pak
usr/lib/das-keyboard-q/locales/ar.pak
usr/lib/das-keyboard-q/locales/bg.pak
usr/lib/das-keyboard-q/locales/bn.pak
usr/lib/das-keyboard-q/locales/ca.pak
usr/lib/das-keyboard-q/locales/cs.pak
usr/lib/das-keyboard-q/locales/da.pak
usr/lib/das-keyboard-q/locales/de.pak
usr/lib/das-keyboard-q/locales/el.pak
usr/lib/das-keyboard-q/locales/en-GB.pak
usr/lib/das-keyboard-q/locales/en-US.pak
usr/lib/das-keyboard-q/locales/es-419.pak
usr/lib/das-keyboard-q/locales/es.pak
usr/lib/das-keyboard-q/locales/et.pak
usr/lib/das-keyboard-q/locales/fa.pak
usr/lib/das-keyboard-q/locales/fi.pak
usr/lib/das-keyboard-q/locales/fil.pak
usr/lib/das-keyboard-q/locales/fr.pak
usr/lib/das-keyboard-q/locales/gu.pak
usr/lib/das-keyboard-q/locales/he.pak
usr/lib/das-keyboard-q/locales/hi.pak
usr/lib/das-keyboard-q/locales/hr.pak
usr/lib/das-keyboard-q/locales/hu.pak
usr/lib/das-keyboard-q/locales/id.pak
usr/lib/das-keyboard-q/locales/it.pak
usr/lib/das-keyboard-q/locales/ja.pak
usr/lib/das-keyboard-q/locales/kn.pak
usr/lib/das-keyboard-q/locales/ko.pak
usr/lib/das-keyboard-q/locales/lt.pak
usr/lib/das-keyboard-q/locales/lv.pak
usr/lib/das-keyboard-q/locales/ml.pak
usr/lib/das-keyboard-q/locales/mr.pak
usr/lib/das-keyboard-q/locales/ms.pak
usr/lib/das-keyboard-q/locales/nb.pak
usr/lib/das-keyboard-q/locales/nl.pak
usr/lib/das-keyboard-q/locales/pl.pak
usr/lib/das-keyboard-q/locales/pt-BR.pak
usr/lib/das-keyboard-q/locales/pt-PT.pak
usr/lib/das-keyboard-q/locales/ro.pak
usr/lib/das-keyboard-q/locales/ru.pak
usr/lib/das-keyboard-q/locales/sk.pak
usr/lib/das-keyboard-q/locales/sl.pak
usr/lib/das-keyboard-q/locales/sr.pak
usr/lib/das-keyboard-q/locales/sv.pak
usr/lib/das-keyboard-q/locales/sw.pak
usr/lib/das-keyboard-q/locales/ta.pak
usr/lib/das-keyboard-q/locales/te.pak
usr/lib/das-keyboard-q/locales/th.pak
usr/lib/das-keyboard-q/locales/tr.pak
usr/lib/das-keyboard-q/locales/uk.pak
usr/lib/das-keyboard-q/locales/vi.pak
usr/lib/das-keyboard-q/locales/zh-CN.pak
usr/lib/das-keyboard-q/locales/zh-TW.pak
usr/lib/das-keyboard-q/resources.pak
usr/lib/das-keyboard-q/resources/
usr/lib/das-keyboard-q/resources/app.asar
usr/lib/das-keyboard-q/resources/app.asar.unpacked/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/externalLibraries/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/externalLibraries/DKWindowsUserServicesHelper.exe
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/externalLibraries/Das_Keyboard_User_Services.dll
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/16x16.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/5Q-box-back.jpg
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/DK4QPID-home-icon.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/DK5QPID-home-icon.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/DK5QSPID-home-icon.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/X50QPID-home-icon.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/bg.gif
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/dk5-q-blurry.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/gray-q-icon-no-bracket (copy).png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/gray-q-icon-no-bracket.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/gray-q-icon.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/icon.icns
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/ifttt/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/ifttt/DK4QPID.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/ifttt/DK5QPID.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/ifttt/X50QPID.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5Q/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5Q/5Q-da.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5Q/5Q-de-DE.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5Q/5Q-en-UK.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5Q/5Q-en-US-MAC.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5Q/5Q-en-US.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5Q/5Q-es-ES.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5Q/5Q-fr-FR.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5Q/5Q-it-IT.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5Q/5Q-ru.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5QS/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5QS/5QS-da.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5QS/5QS-de-DE.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5QS/5QS-en-UK.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5QS/5QS-en-US-MAC.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5QS/5QS-en-US.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5QS/5QS-es-ES.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5QS/5QS-fr-FR.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5QS/5QS-it-IT.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/5QS/5QS-ru.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/DK4Q/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/DK4Q/DK4Q-de-DE.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/DK4Q/DK4Q-en-UK.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/DK4Q/DK4Q-en-US-MAC.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/DK4Q/DK4Q-en-US.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/DK4Q/DK4Q-no.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50-en-UK-drop.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50-en-UK.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/X50-DEFAULT.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/X50-da.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/X50-de-DE.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/X50-en-UK.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/X50-en-US-MAC.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/X50-en-US.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/X50-es-ES.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/X50-fr-FR.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/X50-it-IT.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/layouts/X50/X50-ru.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/load.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/logo-big.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/logo.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/q-icon-no-bracket.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/q-icon.gif
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/q-icon.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/q-iconbackup.gif
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/images/stars.png
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/scripts/
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/scripts/linux.sh
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/scripts/mac.scpt
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/scripts/mac.sh
usr/lib/das-keyboard-q/resources/app.asar.unpacked/assets/scripts/windows.ps1
usr/lib/das-keyboard-q/snapshot_blob.bin
usr/lib/das-keyboard-q/swiftshader/
usr/lib/das-keyboard-q/swiftshader/libEGL.so
usr/lib/das-keyboard-q/swiftshader/libGLESv2.so
usr/lib/das-keyboard-q/v8_context_snapshot.bin
usr/lib/das-keyboard-q/version
usr/lib/das-keyboard-q/vk_swiftshader_icd.json
usr/share/
usr/share/applications/
usr/share/applications/das-keyboard-q.desktop
usr/share/doc/
usr/share/doc/das-keyboard-q/
usr/share/doc/das-keyboard-q/copyright
usr/share/lintian/
usr/share/lintian/overrides/
usr/share/lintian/overrides/das-keyboard-q
usr/share/pixmaps/
usr/share/pixmaps/das-keyboard-q.png
