# files/lokinet-0.9.9-1/files
%FILES%
usr/
usr/bin/
usr/bin/lokinet
usr/bin/lokinet-bootstrap
usr/bin/lokinet-vpn
usr/include/
usr/include/lokimq/
usr/include/lokimq/address.h
usr/include/lokimq/auth.h
usr/include/lokimq/base32z.h
usr/include/lokimq/base64.h
usr/include/lokimq/batch.h
usr/include/lokimq/bt_serialize.h
usr/include/lokimq/bt_value.h
usr/include/lokimq/connections.h
usr/include/lokimq/hex.h
usr/include/lokimq/lokimq.h
usr/include/lokimq/message.h
usr/include/lokimq/variant.h
usr/include/lokimq/version.h
usr/include/oxenmq/
usr/include/oxenmq/address.h
usr/include/oxenmq/auth.h
usr/include/oxenmq/base32z.h
usr/include/oxenmq/base64.h
usr/include/oxenmq/batch.h
usr/include/oxenmq/bt_serialize.h
usr/include/oxenmq/bt_value.h
usr/include/oxenmq/connections.h
usr/include/oxenmq/hex.h
usr/include/oxenmq/message.h
usr/include/oxenmq/oxenmq.h
usr/include/oxenmq/variant.h
usr/include/oxenmq/version.h
usr/include/oxenmq/zmq.hpp
usr/lib/
usr/lib/liboxenmq.so
usr/lib/liboxenmq.so.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/liblokimq.pc
usr/lib/pkgconfig/liboxenmq.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/lokinet-bootstrap.service
usr/lib/systemd/system/lokinet-default-config.service
usr/lib/systemd/system/lokinet-resume.service
usr/lib/systemd/system/lokinet-vpn@.service
usr/lib/systemd/system/lokinet.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/lokinet.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/lokinet.conf
usr/share/
usr/share/licenses/
usr/share/licenses/lokinet/
usr/share/licenses/lokinet/LICENSE
usr/share/polkit-1/
usr/share/polkit-1/rules.d/
usr/share/polkit-1/rules.d/lokinet.rules
var/
var/lib/
var/lib/polkit-1/
var/lib/polkit-1/localauthority/
var/lib/polkit-1/localauthority/10-vendor.d/
var/lib/polkit-1/localauthority/10-vendor.d/lokinet.pkla
