# files/podman-docker-git-4.2.0_dev.r15059.g2dcf3067e-1/files
%FILES%
usr/
usr/bin/
usr/bin/docker
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/user/
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/podman-docker.conf
usr/share/
usr/share/licenses/
usr/share/licenses/podman-docker-git/
usr/share/licenses/podman-docker-git/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/docker-attach.1.gz
usr/share/man/man1/docker-auto-update.1.gz
usr/share/man/man1/docker-build.1.gz
usr/share/man/man1/docker-commit.1.gz
usr/share/man/man1/docker-completion.1.gz
usr/share/man/man1/docker-container-checkpoint.1.gz
usr/share/man/man1/docker-container-cleanup.1.gz
usr/share/man/man1/docker-container-clone.1.gz
usr/share/man/man1/docker-container-diff.1.gz
usr/share/man/man1/docker-container-exists.1.gz
usr/share/man/man1/docker-container-inspect.1.gz
usr/share/man/man1/docker-container-prune.1.gz
usr/share/man/man1/docker-container-restore.1.gz
usr/share/man/man1/docker-container-runlabel.1.gz
usr/share/man/man1/docker-container.1.gz
usr/share/man/man1/docker-cp.1.gz
usr/share/man/man1/docker-create.1.gz
usr/share/man/man1/docker-diff.1.gz
usr/share/man/man1/docker-events.1.gz
usr/share/man/man1/docker-exec.1.gz
usr/share/man/man1/docker-export.1.gz
usr/share/man/man1/docker-generate-kube.1.gz
usr/share/man/man1/docker-generate-systemd.1.gz
usr/share/man/man1/docker-generate.1.gz
usr/share/man/man1/docker-healthcheck-run.1.gz
usr/share/man/man1/docker-healthcheck.1.gz
usr/share/man/man1/docker-history.1.gz
usr/share/man/man1/docker-image-diff.1.gz
usr/share/man/man1/docker-image-exists.1.gz
usr/share/man/man1/docker-image-inspect.1.gz
usr/share/man/man1/docker-image-mount.1.gz
usr/share/man/man1/docker-image-prune.1.gz
usr/share/man/man1/docker-image-scp.1.gz
usr/share/man/man1/docker-image-sign.1.gz
usr/share/man/man1/docker-image-tree.1.gz
usr/share/man/man1/docker-image-trust.1.gz
usr/share/man/man1/docker-image-unmount.1.gz
usr/share/man/man1/docker-image.1.gz
usr/share/man/man1/docker-images.1.gz
usr/share/man/man1/docker-import.1.gz
usr/share/man/man1/docker-info.1.gz
usr/share/man/man1/docker-init.1.gz
usr/share/man/man1/docker-inspect.1.gz
usr/share/man/man1/docker-kill.1.gz
usr/share/man/man1/docker-load.1.gz
usr/share/man/man1/docker-login.1.gz
usr/share/man/man1/docker-logout.1.gz
usr/share/man/man1/docker-logs.1.gz
usr/share/man/man1/docker-machine-init.1.gz
usr/share/man/man1/docker-machine-inspect.1.gz
usr/share/man/man1/docker-machine-list.1.gz
usr/share/man/man1/docker-machine-rm.1.gz
usr/share/man/man1/docker-machine-set.1.gz
usr/share/man/man1/docker-machine-ssh.1.gz
usr/share/man/man1/docker-machine-start.1.gz
usr/share/man/man1/docker-machine-stop.1.gz
usr/share/man/man1/docker-machine.1.gz
usr/share/man/man1/docker-manifest-add.1.gz
usr/share/man/man1/docker-manifest-annotate.1.gz
usr/share/man/man1/docker-manifest-create.1.gz
usr/share/man/man1/docker-manifest-exists.1.gz
usr/share/man/man1/docker-manifest-inspect.1.gz
usr/share/man/man1/docker-manifest-push.1.gz
usr/share/man/man1/docker-manifest-remove.1.gz
usr/share/man/man1/docker-manifest-rm.1.gz
usr/share/man/man1/docker-manifest.1.gz
usr/share/man/man1/docker-mount.1.gz
usr/share/man/man1/docker-network-connect.1.gz
usr/share/man/man1/docker-network-create.1.gz
usr/share/man/man1/docker-network-disconnect.1.gz
usr/share/man/man1/docker-network-exists.1.gz
usr/share/man/man1/docker-network-inspect.1.gz
usr/share/man/man1/docker-network-ls.1.gz
usr/share/man/man1/docker-network-prune.1.gz
usr/share/man/man1/docker-network-reload.1.gz
usr/share/man/man1/docker-network-rm.1.gz
usr/share/man/man1/docker-network.1.gz
usr/share/man/man1/docker-pause.1.gz
usr/share/man/man1/docker-play-kube.1.gz
usr/share/man/man1/docker-play.1.gz
usr/share/man/man1/docker-pod-create.1.gz
usr/share/man/man1/docker-pod-exists.1.gz
usr/share/man/man1/docker-pod-inspect.1.gz
usr/share/man/man1/docker-pod-kill.1.gz
usr/share/man/man1/docker-pod-logs.1.gz
usr/share/man/man1/docker-pod-pause.1.gz
usr/share/man/man1/docker-pod-prune.1.gz
usr/share/man/man1/docker-pod-ps.1.gz
usr/share/man/man1/docker-pod-restart.1.gz
usr/share/man/man1/docker-pod-rm.1.gz
usr/share/man/man1/docker-pod-start.1.gz
usr/share/man/man1/docker-pod-stats.1.gz
usr/share/man/man1/docker-pod-stop.1.gz
usr/share/man/man1/docker-pod-top.1.gz
usr/share/man/man1/docker-pod-unpause.1.gz
usr/share/man/man1/docker-pod.1.gz
usr/share/man/man1/docker-port.1.gz
usr/share/man/man1/docker-ps.1.gz
usr/share/man/man1/docker-pull.1.gz
usr/share/man/man1/docker-push.1.gz
usr/share/man/man1/docker-remote.1.gz
usr/share/man/man1/docker-rename.1.gz
usr/share/man/man1/docker-restart.1.gz
usr/share/man/man1/docker-rm.1.gz
usr/share/man/man1/docker-rmi.1.gz
usr/share/man/man1/docker-run.1.gz
usr/share/man/man1/docker-save.1.gz
usr/share/man/man1/docker-search.1.gz
usr/share/man/man1/docker-secret-create.1.gz
usr/share/man/man1/docker-secret-inspect.1.gz
usr/share/man/man1/docker-secret-ls.1.gz
usr/share/man/man1/docker-secret-rm.1.gz
usr/share/man/man1/docker-secret.1.gz
usr/share/man/man1/docker-start.1.gz
usr/share/man/man1/docker-stats.1.gz
usr/share/man/man1/docker-stop.1.gz
usr/share/man/man1/docker-system-connection-add.1.gz
usr/share/man/man1/docker-system-connection-default.1.gz
usr/share/man/man1/docker-system-connection-list.1.gz
usr/share/man/man1/docker-system-connection-remove.1.gz
usr/share/man/man1/docker-system-connection-rename.1.gz
usr/share/man/man1/docker-system-connection.1.gz
usr/share/man/man1/docker-system-df.1.gz
usr/share/man/man1/docker-system-migrate.1.gz
usr/share/man/man1/docker-system-prune.1.gz
usr/share/man/man1/docker-system-renumber.1.gz
usr/share/man/man1/docker-system-reset.1.gz
usr/share/man/man1/docker-system-service.1.gz
usr/share/man/man1/docker-system.1.gz
usr/share/man/man1/docker-tag.1.gz
usr/share/man/man1/docker-top.1.gz
usr/share/man/man1/docker-unmount.1.gz
usr/share/man/man1/docker-unpause.1.gz
usr/share/man/man1/docker-unshare.1.gz
usr/share/man/man1/docker-untag.1.gz
usr/share/man/man1/docker-version.1.gz
usr/share/man/man1/docker-volume-create.1.gz
usr/share/man/man1/docker-volume-exists.1.gz
usr/share/man/man1/docker-volume-export.1.gz
usr/share/man/man1/docker-volume-import.1.gz
usr/share/man/man1/docker-volume-inspect.1.gz
usr/share/man/man1/docker-volume-ls.1.gz
usr/share/man/man1/docker-volume-mount.1.gz
usr/share/man/man1/docker-volume-prune.1.gz
usr/share/man/man1/docker-volume-rm.1.gz
usr/share/man/man1/docker-volume-unmount.1.gz
usr/share/man/man1/docker-volume.1.gz
usr/share/man/man1/docker-wait.1.gz
usr/share/man/man1/docker.1.gz
usr/share/man/man5/
usr/share/man/man5/dockerfile.5.gz
usr/share/man/man5/dockerignore.5.gz
