# files/libcgroup-2.0-1/files
%FILES%
etc/
etc/cgconfig.conf
etc/cgrules.conf
etc/cgsnapshot_blacklist.conf
usr/
usr/bin/
usr/bin/cgclassify
usr/bin/cgclear
usr/bin/cgconfigparser
usr/bin/cgcreate
usr/bin/cgdelete
usr/bin/cgexec
usr/bin/cgget
usr/bin/cgrulesengd
usr/bin/cgset
usr/bin/cgsnapshot
usr/bin/lscgroup
usr/bin/lssubsys
usr/include/
usr/include/libcgroup.h
usr/include/libcgroup/
usr/include/libcgroup/config.h
usr/include/libcgroup/error.h
usr/include/libcgroup/groups.h
usr/include/libcgroup/init.h
usr/include/libcgroup/iterators.h
usr/include/libcgroup/log.h
usr/include/libcgroup/tasks.h
usr/lib/
usr/lib/libcgroup.so
usr/lib/libcgroup.so.1
usr/lib/libcgroup.so.1.0.42
usr/lib/libcgroupfortesting.so
usr/lib/libcgroupfortesting.so.1
usr/lib/libcgroupfortesting.so.1.0.42
usr/lib/libcgset.so
usr/lib/libcgset.so.0
usr/lib/libcgset.so.0.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/libcgroup.pc
usr/lib/security/
usr/lib/security/pam_cgroup.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/cgconfig.service
usr/lib/systemd/system/cgrules.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/cgclassify.1.gz
usr/share/man/man1/cgclear.1.gz
usr/share/man/man1/cgcreate.1.gz
usr/share/man/man1/cgdelete.1.gz
usr/share/man/man1/cgexec.1.gz
usr/share/man/man1/cgget.1.gz
usr/share/man/man1/cgset.1.gz
usr/share/man/man1/cgsnapshot.1.gz
usr/share/man/man1/lscgroup.1.gz
usr/share/man/man1/lssubsys.1.gz
usr/share/man/man5/
usr/share/man/man5/cgconfig.conf.5.gz
usr/share/man/man5/cgred.conf.5.gz
usr/share/man/man5/cgrules.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/cgconfigparser.8.gz
usr/share/man/man8/cgrulesengd.8.gz
