# files/keymapper-git-r68.d84fd42-1/files
%FILES%
etc/
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/keymapper.desktop
usr/
usr/bin/
usr/bin/keymapper
usr/bin/keymapperd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/keymapperd.service
