# files/croc-git-9.5.3.r3.g35652e6-1/files
%FILES%
usr/
usr/bin/
usr/bin/croc
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/croc.service
usr/share/
usr/share/licenses/
usr/share/licenses/croc/
usr/share/licenses/croc/LICENSE
