# files/preload-0.6.4-7/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/preload
etc/logrotate.d/
etc/logrotate.d/preload
etc/preload.conf
usr/
usr/bin/
usr/bin/preload
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/preload.service
usr/share/
usr/share/doc/
usr/share/doc/preload-0.6.4/
usr/share/doc/preload-0.6.4/index.txt
usr/share/doc/preload-0.6.4/proposal.txt
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/preload.8.gz
var/
var/lib/
var/lib/preload/
var/log/
