# files/pure-ftpd-1.0.49-5/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/pure-ftpd
etc/pure-ftpd/
etc/pure-ftpd/pure-ftpd.conf
etc/pure-ftpd/welcome.msg
usr/
usr/bin/
usr/bin/pure-authd
usr/bin/pure-certd
usr/bin/pure-ftpd
usr/bin/pure-ftpwho
usr/bin/pure-mrtginfo
usr/bin/pure-pw
usr/bin/pure-pwconvert
usr/bin/pure-quotacheck
usr/bin/pure-statsdecode
usr/bin/pure-uploadscript
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pure-ftpd.service
usr/share/
usr/share/doc/
usr/share/doc/pure-ftpd/
usr/share/doc/pure-ftpd/ChangeLog
usr/share/doc/pure-ftpd/FAQ
usr/share/doc/pure-ftpd/README
usr/share/doc/pure-ftpd/README.Authentication-Modules
usr/share/doc/pure-ftpd/README.Configuration-File
usr/share/doc/pure-ftpd/README.Donations
usr/share/doc/pure-ftpd/README.LDAP
usr/share/doc/pure-ftpd/README.MacOS-X
usr/share/doc/pure-ftpd/README.MySQL
usr/share/doc/pure-ftpd/README.PGSQL
usr/share/doc/pure-ftpd/README.TLS
usr/share/doc/pure-ftpd/README.Virtual-Users
usr/share/doc/pure-ftpd/README.Windows
usr/share/doc/pure-ftpd/pure-ftpd.conf
usr/share/doc/pure-ftpd/pureftpd-ldap.conf
usr/share/doc/pure-ftpd/pureftpd-mysql.conf
usr/share/doc/pure-ftpd/pureftpd-pgsql.conf
usr/share/doc/pure-ftpd/pureftpd.schema
usr/share/licenses/
usr/share/licenses/pure-ftpd/
usr/share/licenses/pure-ftpd/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/pure-authd.8.gz
usr/share/man/man8/pure-certd.8.gz
usr/share/man/man8/pure-ftpd.8.gz
usr/share/man/man8/pure-ftpwho.8.gz
usr/share/man/man8/pure-mrtginfo.8.gz
usr/share/man/man8/pure-pw.8.gz
usr/share/man/man8/pure-pwconvert.8.gz
usr/share/man/man8/pure-quotacheck.8.gz
usr/share/man/man8/pure-statsdecode.8.gz
usr/share/man/man8/pure-uploadscript.8.gz
