# files/amdguid-wayland-bin-1.0.11-2/files
%FILES%
usr/
usr/bin/
usr/bin/amdgui-helper
usr/bin/amdguid
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/amdgui-helper.service
