# files/rustdesk-1.1.6-1/files
%FILES%
usr/
usr/bin/
usr/bin/libsciter-gtk.so
usr/bin/rustdesk
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rustdesk.service
usr/share/
usr/share/applications/
usr/share/applications/rustdesk.desktop
usr/share/pixmaps/
usr/share/pixmaps/rustdesk.svg
