# files/system76-power-1.1.22-1/files
%FILES%
usr/
usr/bin/
usr/bin/system76-power
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/com.system76.PowerDaemon.service
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/com.system76.PowerDaemon.xml
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/com.system76.PowerDaemon.conf
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/com.system76.PowerDaemon.policy
