# files/libreswan-4.6-1/files
%FILES%
etc/
etc/ipsec.conf
etc/ipsec.d/
etc/ipsec.d/policies/
etc/ipsec.d/policies/block
etc/ipsec.d/policies/clear
etc/ipsec.d/policies/clear-or-private
etc/ipsec.d/policies/portexcludes.conf
etc/ipsec.d/policies/private
etc/ipsec.d/policies/private-or-clear
etc/ipsec.secrets
etc/logrotate.d/
etc/logrotate.d/libreswan
etc/pam.d/
etc/pam.d/pluto
usr/
usr/bin/
usr/bin/ipsec
usr/lib/
usr/lib/ipsec/
usr/lib/ipsec/_import_crl
usr/lib/ipsec/_plutorun
usr/lib/ipsec/_secretcensor
usr/lib/ipsec/_stackmanager
usr/lib/ipsec/_unbound-hook
usr/lib/ipsec/_updown
usr/lib/ipsec/_updown.xfrm
usr/lib/ipsec/addconn
usr/lib/ipsec/algparse
usr/lib/ipsec/auto
usr/lib/ipsec/barf
usr/lib/ipsec/cavp
usr/lib/ipsec/dncheck
usr/lib/ipsec/ecdsasigkey
usr/lib/ipsec/enumcheck
usr/lib/ipsec/hunkcheck
usr/lib/ipsec/ipcheck
usr/lib/ipsec/jambufcheck
usr/lib/ipsec/keyidcheck
usr/lib/ipsec/letsencrypt
usr/lib/ipsec/look
usr/lib/ipsec/newhostkey
usr/lib/ipsec/pluto
usr/lib/ipsec/readwriteconf
usr/lib/ipsec/rsasigkey
usr/lib/ipsec/setup
usr/lib/ipsec/show
usr/lib/ipsec/showhostkey
usr/lib/ipsec/timecheck
usr/lib/ipsec/verify
usr/lib/ipsec/whack
usr/lib/sysctl.d/
usr/lib/sysctl.d/libreswan-icmp-redirects.conf
usr/lib/sysctl.d/libreswan-icmp-redirects.conf.revert
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ipsec.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/libreswan.conf
usr/share/
usr/share/doc/
usr/share/doc/libreswan/
usr/share/doc/libreswan/ipsec.conf-sample
usr/share/doc/libreswan/ipsec.secrets-sample
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/ipsec.conf.5.gz
usr/share/man/man5/ipsec.secrets.5.gz
usr/share/man/man8/
usr/share/man/man8/ipsec.8.gz
usr/share/man/man8/ipsec__import_crl.8.gz
usr/share/man/man8/ipsec__plutorun.8.gz
usr/share/man/man8/ipsec__secretcensor.8.gz
usr/share/man/man8/ipsec__stackmanager.8.gz
usr/share/man/man8/ipsec__unbound-hook.8.gz
usr/share/man/man8/ipsec__updown.8.gz
usr/share/man/man8/ipsec__updown.xfrm.8.gz
usr/share/man/man8/ipsec_addconn.8.gz
usr/share/man/man8/ipsec_auto.8.gz
usr/share/man/man8/ipsec_barf.8.gz
usr/share/man/man8/ipsec_checknss.8.gz
usr/share/man/man8/ipsec_ecdsasigkey.8.gz
usr/share/man/man8/ipsec_import.8.gz
usr/share/man/man8/ipsec_initnss.8.gz
usr/share/man/man8/ipsec_letsencrypt.8.gz
usr/share/man/man8/ipsec_look.8.gz
usr/share/man/man8/ipsec_newhostkey.8.gz
usr/share/man/man8/ipsec_pluto.8.gz
usr/share/man/man8/ipsec_readwriteconf.8.gz
usr/share/man/man8/ipsec_rsasigkey.8.gz
usr/share/man/man8/ipsec_setup.8.gz
usr/share/man/man8/ipsec_show.8.gz
usr/share/man/man8/ipsec_showhostkey.8.gz
usr/share/man/man8/ipsec_verify.8.gz
usr/share/man/man8/ipsec_whack.8.gz
usr/share/man/man8/pluto.8.gz
var/
var/lib/
var/lib/ipsec/
var/lib/ipsec/nss/
