# files/minetest-server-git-1:5.6.0.r78.g22cbc0580-1/files
%FILES%
etc/
etc/minetest/
usr/
usr/bin/
usr/bin/minetestserver
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/minetest@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/minetest-server.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/minetest-server.conf
usr/share/
usr/share/man/
usr/share/man/man6/
usr/share/man/man6/minetestserver.6.gz
