# files/mozillavpn-2.9.0-2/files
%FILES%
usr/
usr/bin/
usr/bin/mozillavpn
usr/etc/
usr/etc/chromium/
usr/etc/chromium/native-messaging-hosts/
usr/etc/chromium/native-messaging-hosts/mozillavpn.json
usr/etc/opt/
usr/etc/opt/chrome/
usr/etc/opt/chrome/native-messaging-hosts/
usr/etc/opt/chrome/native-messaging-hosts/mozillavpn.json
usr/etc/xdg/
usr/etc/xdg/autostart/
usr/etc/xdg/autostart/mozillavpn-startup.desktop
usr/lib/
usr/lib/mozilla/
usr/lib/mozilla/native-messaging-hosts/
usr/lib/mozilla/native-messaging-hosts/mozillavpn.json
usr/lib/mozillavpn/
usr/lib/mozillavpn/mozillavpnnp
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mozillavpn.service
usr/share/
usr/share/applications/
usr/share/applications/mozillavpn.desktop
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.mozilla.vpn.dbus.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.mozilla.vpn.conf
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/apps/
usr/share/icons/hicolor/128x128/apps/mozillavpn.png
usr/share/icons/hicolor/16x16/
usr/share/icons/hicolor/16x16/apps/
usr/share/icons/hicolor/16x16/apps/mozillavpn.png
usr/share/icons/hicolor/32x32/
usr/share/icons/hicolor/32x32/apps/
usr/share/icons/hicolor/32x32/apps/mozillavpn.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/mozillavpn.png
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/apps/
usr/share/icons/hicolor/64x64/apps/mozillavpn.png
usr/share/mozillavpn/
usr/share/mozillavpn/addons/
usr/share/mozillavpn/addons/guide_01_how_to_vpn.rcc
usr/share/mozillavpn/addons/guide_02_is_my_vpn_working_correctly.rcc
usr/share/mozillavpn/addons/guide_03_adding_and_removing_devices.rcc
usr/share/mozillavpn/addons/guide_04_connecting_external_devices.rcc
usr/share/mozillavpn/addons/tutorial_01_get_started.rcc
usr/share/mozillavpn/addons/tutorial_02_connect_on_startup.rcc
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/org.mozilla.vpn.policy
