# files/ananicy-cpp-0.8.1-1/files
%FILES%
etc/
etc/ananicy.d/
usr/
usr/bin/
usr/bin/ananicy-cpp
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ananicy-cpp.service
