# files/greetd-0.8.0-1/files
%FILES%
etc/
etc/greetd/
etc/greetd/config.toml
etc/pam.d/
etc/pam.d/greetd
usr/
usr/bin/
usr/bin/agreety
usr/bin/greetd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/greetd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/greetd.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/agreety.1.gz
usr/share/man/man1/greetd.1.gz
usr/share/man/man5/
usr/share/man/man5/greetd.5.gz
usr/share/man/man7/
usr/share/man/man7/greetd-ipc.7.gz
