# files/yp-tools-4.2.3-6/files
%FILES%
etc/
etc/nicknames
etc/nisdomainname
usr/
usr/bin/
usr/bin/domainname
usr/bin/nisdomainname
usr/bin/yp_dump_binding
usr/bin/ypcat
usr/bin/ypchfn
usr/bin/ypchsh
usr/bin/ypdomainname
usr/bin/ypmatch
usr/bin/yppasswd
usr/bin/yppoll
usr/bin/ypset
usr/bin/yptest
usr/bin/ypwhich
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/domainname.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/yp-tools.conf
usr/share/
usr/share/locale/
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/yp-tools.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/yp-tools.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/ypcat.1.gz
usr/share/man/man1/ypchfn.1.gz
usr/share/man/man1/ypchsh.1.gz
usr/share/man/man1/ypmatch.1.gz
usr/share/man/man1/yppasswd.1.gz
usr/share/man/man1/ypwhich.1.gz
usr/share/man/man5/
usr/share/man/man5/nicknames.5.gz
usr/share/man/man8/
usr/share/man/man8/domainname.8.gz
usr/share/man/man8/nisdomainname.8.gz
usr/share/man/man8/yp_dump_binding.8.gz
usr/share/man/man8/ypdomainname.8.gz
usr/share/man/man8/yppoll.8.gz
usr/share/man/man8/ypset.8.gz
