# files/syncthing-relaysrv-git-1.15.0.rc.5.r23.g924b96856-1/files
%FILES%
usr/
usr/bin/
usr/bin/syncthing-relaysrv
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/syncthing-relaysrv.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/syncthing-relaysrv.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/syncthing-relaysrv.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/syncthing-relaysrv.1.gz
