# files/jellyfin-server-git-10.8.0.alpha5.r561.g1b020b8abb-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/jellyfin
usr/
usr/lib/
usr/lib/jellyfin/
usr/lib/jellyfin/BDInfo.dll
usr/lib/jellyfin/BlurHashSharp.SkiaSharp.dll
usr/lib/jellyfin/BlurHashSharp.dll
usr/lib/jellyfin/CommandLine.dll
usr/lib/jellyfin/DiscUtils.Core.dll
usr/lib/jellyfin/DiscUtils.Iso9660.dll
usr/lib/jellyfin/DiscUtils.Streams.dll
usr/lib/jellyfin/DiscUtils.Udf.dll
usr/lib/jellyfin/DotNet.Glob.dll
usr/lib/jellyfin/DvdLib.dll
usr/lib/jellyfin/DvdLib.pdb
usr/lib/jellyfin/DvdLib.xml
usr/lib/jellyfin/Emby.Dlna.dll
usr/lib/jellyfin/Emby.Dlna.pdb
usr/lib/jellyfin/Emby.Dlna.xml
usr/lib/jellyfin/Emby.Drawing.dll
usr/lib/jellyfin/Emby.Drawing.pdb
usr/lib/jellyfin/Emby.Drawing.xml
usr/lib/jellyfin/Emby.Naming.dll
usr/lib/jellyfin/Emby.Naming.pdb
usr/lib/jellyfin/Emby.Naming.xml
usr/lib/jellyfin/Emby.Notifications.dll
usr/lib/jellyfin/Emby.Notifications.pdb
usr/lib/jellyfin/Emby.Notifications.xml
usr/lib/jellyfin/Emby.Photos.dll
usr/lib/jellyfin/Emby.Photos.pdb
usr/lib/jellyfin/Emby.Photos.xml
usr/lib/jellyfin/Emby.Server.Implementations.dll
usr/lib/jellyfin/Emby.Server.Implementations.pdb
usr/lib/jellyfin/Emby.Server.Implementations.xml
usr/lib/jellyfin/Jellyfin.Api.dll
usr/lib/jellyfin/Jellyfin.Api.pdb
usr/lib/jellyfin/Jellyfin.Api.xml
usr/lib/jellyfin/Jellyfin.Data.dll
usr/lib/jellyfin/Jellyfin.Data.pdb
usr/lib/jellyfin/Jellyfin.Data.xml
usr/lib/jellyfin/Jellyfin.Drawing.Skia.dll
usr/lib/jellyfin/Jellyfin.Drawing.Skia.pdb
usr/lib/jellyfin/Jellyfin.Drawing.Skia.xml
usr/lib/jellyfin/Jellyfin.Extensions.dll
usr/lib/jellyfin/Jellyfin.Extensions.pdb
usr/lib/jellyfin/Jellyfin.Extensions.xml
usr/lib/jellyfin/Jellyfin.MediaEncoding.Hls.dll
usr/lib/jellyfin/Jellyfin.MediaEncoding.Hls.pdb
usr/lib/jellyfin/Jellyfin.MediaEncoding.Hls.xml
usr/lib/jellyfin/Jellyfin.MediaEncoding.Keyframes.dll
usr/lib/jellyfin/Jellyfin.MediaEncoding.Keyframes.pdb
usr/lib/jellyfin/Jellyfin.MediaEncoding.Keyframes.xml
usr/lib/jellyfin/Jellyfin.Networking.dll
usr/lib/jellyfin/Jellyfin.Networking.pdb
usr/lib/jellyfin/Jellyfin.Networking.xml
usr/lib/jellyfin/Jellyfin.Server.Implementations.dll
usr/lib/jellyfin/Jellyfin.Server.Implementations.pdb
usr/lib/jellyfin/Jellyfin.Server.Implementations.xml
usr/lib/jellyfin/Jellyfin.XmlTv.dll
usr/lib/jellyfin/MediaBrowser.Common.dll
usr/lib/jellyfin/MediaBrowser.Common.pdb
usr/lib/jellyfin/MediaBrowser.Common.xml
usr/lib/jellyfin/MediaBrowser.Controller.dll
usr/lib/jellyfin/MediaBrowser.Controller.pdb
usr/lib/jellyfin/MediaBrowser.Controller.xml
usr/lib/jellyfin/MediaBrowser.LocalMetadata.dll
usr/lib/jellyfin/MediaBrowser.LocalMetadata.pdb
usr/lib/jellyfin/MediaBrowser.LocalMetadata.xml
usr/lib/jellyfin/MediaBrowser.MediaEncoding.dll
usr/lib/jellyfin/MediaBrowser.MediaEncoding.pdb
usr/lib/jellyfin/MediaBrowser.MediaEncoding.xml
usr/lib/jellyfin/MediaBrowser.Model.dll
usr/lib/jellyfin/MediaBrowser.Model.pdb
usr/lib/jellyfin/MediaBrowser.Model.xml
usr/lib/jellyfin/MediaBrowser.Providers.dll
usr/lib/jellyfin/MediaBrowser.Providers.pdb
usr/lib/jellyfin/MediaBrowser.Providers.xml
usr/lib/jellyfin/MediaBrowser.XbmcMetadata.dll
usr/lib/jellyfin/MediaBrowser.XbmcMetadata.pdb
usr/lib/jellyfin/MediaBrowser.XbmcMetadata.xml
usr/lib/jellyfin/Microsoft.AspNetCore.Authorization.dll
usr/lib/jellyfin/Microsoft.AspNetCore.Metadata.dll
usr/lib/jellyfin/Microsoft.Bcl.AsyncInterfaces.dll
usr/lib/jellyfin/Microsoft.Data.Sqlite.dll
usr/lib/jellyfin/Microsoft.EntityFrameworkCore.Abstractions.dll
usr/lib/jellyfin/Microsoft.EntityFrameworkCore.Relational.dll
usr/lib/jellyfin/Microsoft.EntityFrameworkCore.Sqlite.dll
usr/lib/jellyfin/Microsoft.EntityFrameworkCore.dll
usr/lib/jellyfin/Microsoft.Extensions.Caching.Memory.dll
usr/lib/jellyfin/Microsoft.Extensions.Configuration.EnvironmentVariables.dll
usr/lib/jellyfin/Microsoft.Extensions.DependencyModel.dll
usr/lib/jellyfin/Microsoft.Extensions.Diagnostics.HealthChecks.Abstractions.dll
usr/lib/jellyfin/Microsoft.Extensions.Diagnostics.HealthChecks.EntityFrameworkCore.dll
usr/lib/jellyfin/Microsoft.Extensions.Diagnostics.HealthChecks.dll
usr/lib/jellyfin/Microsoft.Extensions.Logging.Abstractions.dll
usr/lib/jellyfin/Microsoft.OpenApi.dll
usr/lib/jellyfin/Microsoft.Win32.SystemEvents.dll
usr/lib/jellyfin/Mono.Nat.dll
usr/lib/jellyfin/NEbml.Core.dll
usr/lib/jellyfin/Newtonsoft.Json.dll
usr/lib/jellyfin/PlaylistsNET.dll
usr/lib/jellyfin/PriorityQueue.dll
usr/lib/jellyfin/Prometheus.AspNetCore.dll
usr/lib/jellyfin/Prometheus.NetCore.dll
usr/lib/jellyfin/Prometheus.NetStandard.dll
usr/lib/jellyfin/RSSDP.dll
usr/lib/jellyfin/RSSDP.pdb
usr/lib/jellyfin/Resources/
usr/lib/jellyfin/Resources/Configuration/
usr/lib/jellyfin/Resources/Configuration/logging.json
usr/lib/jellyfin/SQLitePCL.pretty.dll
usr/lib/jellyfin/SQLitePCLRaw.batteries_v2.dll
usr/lib/jellyfin/SQLitePCLRaw.core.dll
usr/lib/jellyfin/SQLitePCLRaw.provider.e_sqlite3.dll
usr/lib/jellyfin/Serilog.AspNetCore.dll
usr/lib/jellyfin/Serilog.Enrichers.Thread.dll
usr/lib/jellyfin/Serilog.Extensions.Hosting.dll
usr/lib/jellyfin/Serilog.Extensions.Logging.dll
usr/lib/jellyfin/Serilog.Formatting.Compact.dll
usr/lib/jellyfin/Serilog.Settings.Configuration.dll
usr/lib/jellyfin/Serilog.Sinks.Async.dll
usr/lib/jellyfin/Serilog.Sinks.Console.dll
usr/lib/jellyfin/Serilog.Sinks.Debug.dll
usr/lib/jellyfin/Serilog.Sinks.File.dll
usr/lib/jellyfin/Serilog.Sinks.Graylog.Core.dll
usr/lib/jellyfin/Serilog.Sinks.Graylog.dll
usr/lib/jellyfin/Serilog.dll
usr/lib/jellyfin/SharpCompress.dll
usr/lib/jellyfin/SkiaSharp.dll
usr/lib/jellyfin/Swashbuckle.AspNetCore.ReDoc.dll
usr/lib/jellyfin/Swashbuckle.AspNetCore.Swagger.dll
usr/lib/jellyfin/Swashbuckle.AspNetCore.SwaggerGen.dll
usr/lib/jellyfin/Swashbuckle.AspNetCore.SwaggerUI.dll
usr/lib/jellyfin/System.Drawing.Common.dll
usr/lib/jellyfin/System.Linq.Async.dll
usr/lib/jellyfin/System.Text.Json.dll
usr/lib/jellyfin/TMDbLib.dll
usr/lib/jellyfin/TagLibSharp.dll
usr/lib/jellyfin/UtfUnknown.dll
usr/lib/jellyfin/jellyfin
usr/lib/jellyfin/jellyfin.deps.json
usr/lib/jellyfin/jellyfin.dll
usr/lib/jellyfin/jellyfin.pdb
usr/lib/jellyfin/jellyfin.runtimeconfig.json
usr/lib/jellyfin/jellyfin.xml
usr/lib/jellyfin/libse.dll
usr/lib/jellyfin/prometheus-net.DotNetRuntime.dll
usr/lib/jellyfin/runtimes/
usr/lib/jellyfin/runtimes/linux-arm/
usr/lib/jellyfin/runtimes/linux-arm/native/
usr/lib/jellyfin/runtimes/linux-arm/native/libSkiaSharp.so
usr/lib/jellyfin/runtimes/linux-arm/native/libe_sqlite3.so
usr/lib/jellyfin/runtimes/linux-arm64/
usr/lib/jellyfin/runtimes/linux-arm64/native/
usr/lib/jellyfin/runtimes/linux-arm64/native/libSkiaSharp.so
usr/lib/jellyfin/runtimes/linux-arm64/native/libe_sqlite3.so
usr/lib/jellyfin/runtimes/linux-armel/
usr/lib/jellyfin/runtimes/linux-armel/native/
usr/lib/jellyfin/runtimes/linux-armel/native/libe_sqlite3.so
usr/lib/jellyfin/runtimes/linux-mips64/
usr/lib/jellyfin/runtimes/linux-mips64/native/
usr/lib/jellyfin/runtimes/linux-mips64/native/libe_sqlite3.so
usr/lib/jellyfin/runtimes/linux-musl-x64/
usr/lib/jellyfin/runtimes/linux-musl-x64/native/
usr/lib/jellyfin/runtimes/linux-musl-x64/native/libSkiaSharp.so
usr/lib/jellyfin/runtimes/linux-musl-x64/native/libe_sqlite3.so
usr/lib/jellyfin/runtimes/linux-s390x/
usr/lib/jellyfin/runtimes/linux-s390x/native/
usr/lib/jellyfin/runtimes/linux-s390x/native/libe_sqlite3.so
usr/lib/jellyfin/runtimes/linux-x64/
usr/lib/jellyfin/runtimes/linux-x64/native/
usr/lib/jellyfin/runtimes/linux-x64/native/libSkiaSharp.so
usr/lib/jellyfin/runtimes/linux-x64/native/libe_sqlite3.so
usr/lib/jellyfin/runtimes/linux-x86/
usr/lib/jellyfin/runtimes/linux-x86/native/
usr/lib/jellyfin/runtimes/linux-x86/native/libe_sqlite3.so
usr/lib/jellyfin/runtimes/unix/
usr/lib/jellyfin/runtimes/unix/lib/
usr/lib/jellyfin/runtimes/unix/lib/netcoreapp3.0/
usr/lib/jellyfin/runtimes/unix/lib/netcoreapp3.0/System.Drawing.Common.dll
usr/lib/jellyfin/web.config
usr/lib/jellyfin/wwwroot/
usr/lib/jellyfin/wwwroot/api-docs/
usr/lib/jellyfin/wwwroot/api-docs/banner-dark.svg
usr/lib/jellyfin/wwwroot/api-docs/redoc/
usr/lib/jellyfin/wwwroot/api-docs/redoc/custom.css
usr/lib/jellyfin/wwwroot/api-docs/swagger/
usr/lib/jellyfin/wwwroot/api-docs/swagger/custom.css
usr/lib/jellyfin/zlib.net.dll
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/jellyfin.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/jellyfin.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/jellyfin.conf
