# files/mongodb-bin-5.0.9-1/files
%FILES%
etc/
etc/mongodb.conf
usr/
usr/bin/
usr/bin/mongo
usr/bin/mongod
usr/bin/mongos
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mongodb.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/mongodb.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/mongodb.conf
usr/share/
usr/share/licenses/
usr/share/licenses/mongodb-bin/
usr/share/licenses/mongodb-bin/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/mongo.1.gz
usr/share/man/man1/mongod.1.gz
usr/share/man/man1/mongos.1.gz
