# files/foo2zjs-nightly-20201127-1.1/files
%FILES%
usr/
usr/bin/
usr/bin/arm2hpdl
usr/bin/command2foo2lava-pjl
usr/bin/ddstdecode
usr/bin/foo2ddst
usr/bin/foo2ddst-wrapper
usr/bin/foo2hbpl2
usr/bin/foo2hbpl2-wrapper
usr/bin/foo2hiperc
usr/bin/foo2hiperc-wrapper
usr/bin/foo2hp
usr/bin/foo2hp2600-wrapper
usr/bin/foo2lava
usr/bin/foo2lava-wrapper
usr/bin/foo2oak
usr/bin/foo2oak-wrapper
usr/bin/foo2qpdl
usr/bin/foo2qpdl-wrapper
usr/bin/foo2slx
usr/bin/foo2slx-wrapper
usr/bin/foo2xqx
usr/bin/foo2xqx-wrapper
usr/bin/foo2zjs
usr/bin/foo2zjs-icc2ps
usr/bin/foo2zjs-loadfw
usr/bin/foo2zjs-pstops
usr/bin/foo2zjs-wrapper
usr/bin/gipddecode
usr/bin/hbpldecode
usr/bin/hipercdecode
usr/bin/lavadecode
usr/bin/oakdecode
usr/bin/opldecode
usr/bin/printer-profile
usr/bin/qpdldecode
usr/bin/slxdecode
usr/bin/xqxdecode
usr/bin/zjsdecode
usr/lib/
usr/lib/cups/
usr/lib/cups/filter/
usr/lib/cups/filter/command2foo2lava-pjl
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/foo2zjs-loadfw@.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/71-hplj10xx.rules
usr/share/
usr/share/cups/
usr/share/cups/model/
usr/share/cups/model/Dell-1355.ppd.gz
usr/share/cups/model/Dell-C1765.ppd.gz
usr/share/cups/model/Epson-AcuLaser_CX17NF.ppd.gz
usr/share/cups/model/Epson-AcuLaser_M1400.ppd.gz
usr/share/cups/model/Fuji_Xerox-DocuPrint_CM205.ppd.gz
usr/share/cups/model/Fuji_Xerox-DocuPrint_CM215.ppd.gz
usr/share/cups/model/Fuji_Xerox-DocuPrint_M215.ppd.gz
usr/share/cups/model/Fuji_Xerox-DocuPrint_P205.ppd.gz
usr/share/cups/model/Generic-OAKT_Printer.ppd.gz
usr/share/cups/model/Generic-ZjStream_Printer.ppd.gz
usr/share/cups/model/HP-Color_LaserJet_1500.ppd.gz
usr/share/cups/model/HP-Color_LaserJet_1600.ppd.gz
usr/share/cups/model/HP-Color_LaserJet_2600n.ppd.gz
usr/share/cups/model/HP-Color_LaserJet_CP1215.ppd.gz
usr/share/cups/model/HP-LaserJet_1000.ppd.gz
usr/share/cups/model/HP-LaserJet_1005.ppd.gz
usr/share/cups/model/HP-LaserJet_1018.ppd.gz
usr/share/cups/model/HP-LaserJet_1020.ppd.gz
usr/share/cups/model/HP-LaserJet_1022.ppd.gz
usr/share/cups/model/HP-LaserJet_1022n.ppd.gz
usr/share/cups/model/HP-LaserJet_1022nw.ppd.gz
usr/share/cups/model/HP-LaserJet_M1005_MFP.ppd.gz
usr/share/cups/model/HP-LaserJet_M1120_MFP.ppd.gz
usr/share/cups/model/HP-LaserJet_M1319_MFP.ppd.gz
usr/share/cups/model/HP-LaserJet_P1005.ppd.gz
usr/share/cups/model/HP-LaserJet_P1006.ppd.gz
usr/share/cups/model/HP-LaserJet_P1007.ppd.gz
usr/share/cups/model/HP-LaserJet_P1008.ppd.gz
usr/share/cups/model/HP-LaserJet_P1505.ppd.gz
usr/share/cups/model/HP-LaserJet_P1505n.ppd.gz
usr/share/cups/model/HP-LaserJet_P2014.ppd.gz
usr/share/cups/model/HP-LaserJet_P2014n.ppd.gz
usr/share/cups/model/HP-LaserJet_P2035.ppd.gz
usr/share/cups/model/HP-LaserJet_P2035n.ppd.gz
usr/share/cups/model/HP-LaserJet_Pro_CP1025nw.ppd.gz
usr/share/cups/model/HP-LaserJet_Pro_M1132s_MFP.ppd.gz
usr/share/cups/model/HP-LaserJet_Pro_M1212nf_MFP.ppd.gz
usr/share/cups/model/HP-LaserJet_Pro_M12a.ppd.gz
usr/share/cups/model/HP-LaserJet_Pro_M12w.ppd.gz
usr/share/cups/model/HP-LaserJet_Pro_P1102.ppd.gz
usr/share/cups/model/HP-LaserJet_Pro_P1102w.ppd.gz
usr/share/cups/model/HP-LaserJet_Pro_P1566.ppd.gz
usr/share/cups/model/HP-LaserJet_Pro_P1606dn.ppd.gz
usr/share/cups/model/KONICA_MINOLTA-magicolor_1600W.ppd.gz
usr/share/cups/model/KONICA_MINOLTA-magicolor_1680MF.ppd.gz
usr/share/cups/model/KONICA_MINOLTA-magicolor_1690MF.ppd.gz
usr/share/cups/model/KONICA_MINOLTA-magicolor_2430_DL.ppd.gz
usr/share/cups/model/KONICA_MINOLTA-magicolor_2480_MF.ppd.gz
usr/share/cups/model/KONICA_MINOLTA-magicolor_2490_MF.ppd.gz
usr/share/cups/model/KONICA_MINOLTA-magicolor_2530_DL.ppd.gz
usr/share/cups/model/KONICA_MINOLTA-magicolor_4690MF.ppd.gz
usr/share/cups/model/Kyocera-KM-1635.ppd.gz
usr/share/cups/model/Kyocera-KM-2035.ppd.gz
usr/share/cups/model/Lexmark-C500.ppd.gz
usr/share/cups/model/Minolta-Color_PageWorks_Pro_L.ppd.gz
usr/share/cups/model/Minolta-magicolor_2200_DL.ppd.gz
usr/share/cups/model/Minolta-magicolor_2300_DL.ppd.gz
usr/share/cups/model/Minolta-magicolor_2430_DL.ppd.gz
usr/share/cups/model/Oki-C110.ppd.gz
usr/share/cups/model/Oki-C301dn.ppd.gz
usr/share/cups/model/Oki-C3100.ppd.gz
usr/share/cups/model/Oki-C310dn.ppd.gz
usr/share/cups/model/Oki-C3200.ppd.gz
usr/share/cups/model/Oki-C3300.ppd.gz
usr/share/cups/model/Oki-C3400.ppd.gz
usr/share/cups/model/Oki-C3530_MFP.ppd.gz
usr/share/cups/model/Oki-C5100.ppd.gz
usr/share/cups/model/Oki-C511dn.ppd.gz
usr/share/cups/model/Oki-C5200.ppd.gz
usr/share/cups/model/Oki-C5500.ppd.gz
usr/share/cups/model/Oki-C5600.ppd.gz
usr/share/cups/model/Oki-C5650.ppd.gz
usr/share/cups/model/Oki-C5800.ppd.gz
usr/share/cups/model/Oki-C810.ppd.gz
usr/share/cups/model/Olivetti-d-Color_P160W.ppd.gz
usr/share/cups/model/Ricoh-SP_112.ppd.gz
usr/share/cups/model/Ricoh-SP_201Nw.ppd.gz
usr/share/cups/model/Samsung-CLP-300.ppd.gz
usr/share/cups/model/Samsung-CLP-310.ppd.gz
usr/share/cups/model/Samsung-CLP-315.ppd.gz
usr/share/cups/model/Samsung-CLP-325.ppd.gz
usr/share/cups/model/Samsung-CLP-365.ppd.gz
usr/share/cups/model/Samsung-CLP-600.ppd.gz
usr/share/cups/model/Samsung-CLP-610.ppd.gz
usr/share/cups/model/Samsung-CLP-620.ppd.gz
usr/share/cups/model/Samsung-CLX-2160.ppd.gz
usr/share/cups/model/Samsung-CLX-3160.ppd.gz
usr/share/cups/model/Samsung-CLX-3175.ppd.gz
usr/share/cups/model/Samsung-CLX-3185.ppd.gz
usr/share/cups/model/Samsung-ML-1670.ppd.gz
usr/share/cups/model/Samsung-ML-1675.ppd.gz
usr/share/cups/model/Xerox-Phaser_3010.ppd.gz
usr/share/cups/model/Xerox-Phaser_3040.ppd.gz
usr/share/cups/model/Xerox-Phaser_6110.ppd.gz
usr/share/cups/model/Xerox-Phaser_6115MFP.ppd.gz
usr/share/cups/model/Xerox-Phaser_6121MFP.ppd.gz
usr/share/cups/model/Xerox-WorkCentre_3045.ppd.gz
usr/share/cups/model/Xerox-WorkCentre_6015.ppd.gz
usr/share/doc/
usr/share/doc/foo2zjs/
usr/share/doc/foo2zjs/COPYING
usr/share/doc/foo2zjs/ChangeLog
usr/share/doc/foo2zjs/INSTALL
usr/share/doc/foo2zjs/INSTALL.osx
usr/share/doc/foo2zjs/README
usr/share/doc/foo2zjs/manual.pdf
usr/share/foo2hiperc/
usr/share/foo2hiperc/icm/
usr/share/foo2hiperc/icm/C3400-L1-BIN-ICC_0603150.icm
usr/share/foo2hiperc/icm/C3400-L2-BIN-ICC_0604110.icm
usr/share/foo2hiperc/icm/OK3200L1.icm
usr/share/foo2hiperc/icm/OK3200L2.icm
usr/share/foo2hiperc/icm/OK5600L1.icm
usr/share/foo2hiperc/icm/OK5600L2.icm
usr/share/foo2hiperc/icm/OK5600L3.icm
usr/share/foo2hiperc/icm/OKC301H1.icm
usr/share/foo2hiperc/icm/OKC301H2.icm
usr/share/foo2hiperc/icm/OKC301H3.icm
usr/share/foo2hiperc/icm/OKC310H1.icm
usr/share/foo2hiperc/icm/OKC310H2.icm
usr/share/foo2hiperc/icm/OKC310H3.icm
usr/share/foo2hiperc/icm/OKC511H1.icm
usr/share/foo2hiperc/icm/OKC511H2.icm
usr/share/foo2hiperc/icm/OKC511H3.icm
usr/share/foo2hiperc/icm/OKC810H1.icm
usr/share/foo2hiperc/icm/OKC810H2.icm
usr/share/foo2hiperc/icm/OKC810H3.icm
usr/share/foo2hp/
usr/share/foo2hp/icm/
usr/share/foo2hp/icm/hp1215-argyll-0.icm
usr/share/foo2hp/icm/hpclj2600n-0.icm
usr/share/foo2hp/icm/hpclj2600n-1.icm
usr/share/foo2hp/icm/km2430_0.icm
usr/share/foo2hp/icm/km2430_1.icm
usr/share/foo2hp/icm/km2430_2.icm
usr/share/foo2lava/
usr/share/foo2lava/icm/
usr/share/foo2lava/icm/km-1600-rgb-392-bpp1.icm
usr/share/foo2lava/icm/km2530-jconner-d50.icm
usr/share/foo2lava/icm/km2530_0.icm
usr/share/foo2lava/icm/km2530_1.icm
usr/share/foo2lava/icm/km2530_2.icm
usr/share/foo2oak/
usr/share/foo2oak/icm/
usr/share/foo2oak/icm/hpclj2500rgb.icm
usr/share/foo2oak/icm/hpclj2600n-0.icm
usr/share/foo2oak/icm/hpclj2600n-1.icm
usr/share/foo2qpdl/
usr/share/foo2qpdl/crd/
usr/share/foo2qpdl/crd/CLP-300-1200x1200cms2
usr/share/foo2qpdl/crd/CLP-300-1200x600cms2
usr/share/foo2qpdl/crd/CLP-300-600x600cms2
usr/share/foo2qpdl/crd/CLP-300cms
usr/share/foo2qpdl/crd/CLP-600-1200x1200cms2
usr/share/foo2qpdl/crd/CLP-600-1200x600cms2
usr/share/foo2qpdl/crd/CLP-600-600x600cms2
usr/share/foo2qpdl/crd/CLP-600cms
usr/share/foo2qpdl/crd/black-text.ps
usr/share/foo2qpdl/icm/
usr/share/foo2qpdl/icm/samclp300-0.icm
usr/share/foo2qpdl/icm/samclp315-argyll-0.icm
usr/share/foo2slx/
usr/share/foo2slx/icm/
usr/share/foo2slx/icm/lexRPCA2000.icm
usr/share/foo2xqx/
usr/share/foo2xqx/firmware/
usr/share/foo2xqx/firmware/sihpP1005.dl
usr/share/foo2xqx/firmware/sihpP1006.dl
usr/share/foo2xqx/firmware/sihpP1505.dl
usr/share/foo2zjs/
usr/share/foo2zjs/crd/
usr/share/foo2zjs/crd/2300w-1200@150-l250-kx,ucr100,0-per.crd
usr/share/foo2zjs/crd/2300w-1200@150-l250-kx,ucr100,50-per.crd
usr/share/foo2zjs/crd/2300w-1200@150-l250-kx,ucr100,75-per.crd
usr/share/foo2zjs/crd/2300w-1200@150-l250-kx,ucr125,75-per.crd
usr/share/foo2zjs/crd/kh.crd
usr/share/foo2zjs/crd/kx.crd
usr/share/foo2zjs/crd/prolog.ps
usr/share/foo2zjs/crd/screen1200.ps
usr/share/foo2zjs/crd/screen2400.ps
usr/share/foo2zjs/firmware/
usr/share/foo2zjs/firmware/sihp1000.dl
usr/share/foo2zjs/firmware/sihp1005.dl
usr/share/foo2zjs/firmware/sihp1018.dl
usr/share/foo2zjs/firmware/sihp1020.dl
usr/share/foo2zjs/gamma-lookup.ps
usr/share/foo2zjs/gamma.ps
usr/share/foo2zjs/hplj1020_icon.gif
usr/share/foo2zjs/hplj10xx_gui.tcl
usr/share/foo2zjs/icm/
usr/share/foo2zjs/icm/CPWL12W.icm
usr/share/foo2zjs/icm/CPWL24W.icm
usr/share/foo2zjs/icm/CPWL6W.icm
usr/share/foo2zjs/icm/DL2200RGB.icm
usr/share/foo2zjs/icm/DL2312.icm
usr/share/foo2zjs/icm/DL2324.icm
usr/share/foo2zjs/icm/hp-cp1025-rgb-392-bpp1.icm
usr/share/foo2zjs/icm/km2430_0.icm
usr/share/foo2zjs/icm/km2430_1.icm
usr/share/foo2zjs/icm/km2430_2.icm
usr/share/licenses/
usr/share/licenses/foo2zjs-nightly/
usr/share/licenses/foo2zjs-nightly/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/arm2hpdl.1.gz
usr/share/man/man1/ddstdecode.1.gz
usr/share/man/man1/foo2ddst-wrapper.1.gz
usr/share/man/man1/foo2ddst.1.gz
usr/share/man/man1/foo2hbpl2-wrapper.1.gz
usr/share/man/man1/foo2hbpl2.1.gz
usr/share/man/man1/foo2hiperc-wrapper.1.gz
usr/share/man/man1/foo2hiperc.1.gz
usr/share/man/man1/foo2hp.1.gz
usr/share/man/man1/foo2hp2600-wrapper.1.gz
usr/share/man/man1/foo2lava-wrapper.1.gz
usr/share/man/man1/foo2lava.1.gz
usr/share/man/man1/foo2oak-wrapper.1.gz
usr/share/man/man1/foo2oak.1.gz
usr/share/man/man1/foo2qpdl-wrapper.1.gz
usr/share/man/man1/foo2qpdl.1.gz
usr/share/man/man1/foo2slx-wrapper.1.gz
usr/share/man/man1/foo2slx.1.gz
usr/share/man/man1/foo2xqx-wrapper.1.gz
usr/share/man/man1/foo2xqx.1.gz
usr/share/man/man1/foo2zjs-icc2ps.1.gz
usr/share/man/man1/foo2zjs-pstops.1.gz
usr/share/man/man1/foo2zjs-wrapper.1.gz
usr/share/man/man1/foo2zjs.1.gz
usr/share/man/man1/gipddecode.1.gz
usr/share/man/man1/hbpldecode.1.gz
usr/share/man/man1/hipercdecode.1.gz
usr/share/man/man1/lavadecode.1.gz
usr/share/man/man1/oakdecode.1.gz
usr/share/man/man1/opldecode.1.gz
usr/share/man/man1/printer-profile.1.gz
usr/share/man/man1/qpdldecode.1.gz
usr/share/man/man1/slxdecode.1.gz
usr/share/man/man1/usb_printerid.1.gz
usr/share/man/man1/xqxdecode.1.gz
usr/share/man/man1/zjsdecode.1.gz
