# files/ircd-hybrid-8.2.39-1/files
%FILES%
etc/
etc/ircd-hybrid/
etc/ircd-hybrid/ircd.conf
etc/ircd-hybrid/reference.conf
etc/ircd-hybrid/ssl/
run/
run/ircd/
usr/
usr/bin/
usr/bin/ircd
usr/bin/mkpasswd
usr/lib/
usr/lib/ircd-hybrid/
usr/lib/ircd-hybrid/modules/
usr/lib/ircd-hybrid/modules/autoload/
usr/lib/ircd-hybrid/modules/autoload/m_accept.la
usr/lib/ircd-hybrid/modules/autoload/m_accept.so
usr/lib/ircd-hybrid/modules/autoload/m_admin.la
usr/lib/ircd-hybrid/modules/autoload/m_admin.so
usr/lib/ircd-hybrid/modules/autoload/m_away.la
usr/lib/ircd-hybrid/modules/autoload/m_away.so
usr/lib/ircd-hybrid/modules/autoload/m_cap.la
usr/lib/ircd-hybrid/modules/autoload/m_cap.so
usr/lib/ircd-hybrid/modules/autoload/m_capab.la
usr/lib/ircd-hybrid/modules/autoload/m_capab.so
usr/lib/ircd-hybrid/modules/autoload/m_certfp.la
usr/lib/ircd-hybrid/modules/autoload/m_certfp.so
usr/lib/ircd-hybrid/modules/autoload/m_close.la
usr/lib/ircd-hybrid/modules/autoload/m_close.so
usr/lib/ircd-hybrid/modules/autoload/m_connect.la
usr/lib/ircd-hybrid/modules/autoload/m_connect.so
usr/lib/ircd-hybrid/modules/autoload/m_dline.la
usr/lib/ircd-hybrid/modules/autoload/m_dline.so
usr/lib/ircd-hybrid/modules/autoload/m_encap.la
usr/lib/ircd-hybrid/modules/autoload/m_encap.so
usr/lib/ircd-hybrid/modules/autoload/m_eob.la
usr/lib/ircd-hybrid/modules/autoload/m_eob.so
usr/lib/ircd-hybrid/modules/autoload/m_etrace.la
usr/lib/ircd-hybrid/modules/autoload/m_etrace.so
usr/lib/ircd-hybrid/modules/autoload/m_globops.la
usr/lib/ircd-hybrid/modules/autoload/m_globops.so
usr/lib/ircd-hybrid/modules/autoload/m_hash.la
usr/lib/ircd-hybrid/modules/autoload/m_hash.so
usr/lib/ircd-hybrid/modules/autoload/m_help.la
usr/lib/ircd-hybrid/modules/autoload/m_help.so
usr/lib/ircd-hybrid/modules/autoload/m_info.la
usr/lib/ircd-hybrid/modules/autoload/m_info.so
usr/lib/ircd-hybrid/modules/autoload/m_invite.la
usr/lib/ircd-hybrid/modules/autoload/m_invite.so
usr/lib/ircd-hybrid/modules/autoload/m_ison.la
usr/lib/ircd-hybrid/modules/autoload/m_ison.so
usr/lib/ircd-hybrid/modules/autoload/m_kline.la
usr/lib/ircd-hybrid/modules/autoload/m_kline.so
usr/lib/ircd-hybrid/modules/autoload/m_knock.la
usr/lib/ircd-hybrid/modules/autoload/m_knock.so
usr/lib/ircd-hybrid/modules/autoload/m_links.la
usr/lib/ircd-hybrid/modules/autoload/m_links.so
usr/lib/ircd-hybrid/modules/autoload/m_list.la
usr/lib/ircd-hybrid/modules/autoload/m_list.so
usr/lib/ircd-hybrid/modules/autoload/m_locops.la
usr/lib/ircd-hybrid/modules/autoload/m_locops.so
usr/lib/ircd-hybrid/modules/autoload/m_lusers.la
usr/lib/ircd-hybrid/modules/autoload/m_lusers.so
usr/lib/ircd-hybrid/modules/autoload/m_map.la
usr/lib/ircd-hybrid/modules/autoload/m_map.so
usr/lib/ircd-hybrid/modules/autoload/m_metadata.la
usr/lib/ircd-hybrid/modules/autoload/m_metadata.so
usr/lib/ircd-hybrid/modules/autoload/m_module.la
usr/lib/ircd-hybrid/modules/autoload/m_module.so
usr/lib/ircd-hybrid/modules/autoload/m_monitor.la
usr/lib/ircd-hybrid/modules/autoload/m_monitor.so
usr/lib/ircd-hybrid/modules/autoload/m_motd.la
usr/lib/ircd-hybrid/modules/autoload/m_motd.so
usr/lib/ircd-hybrid/modules/autoload/m_names.la
usr/lib/ircd-hybrid/modules/autoload/m_names.so
usr/lib/ircd-hybrid/modules/autoload/m_oper.la
usr/lib/ircd-hybrid/modules/autoload/m_oper.so
usr/lib/ircd-hybrid/modules/autoload/m_pass.la
usr/lib/ircd-hybrid/modules/autoload/m_pass.so
usr/lib/ircd-hybrid/modules/autoload/m_ping.la
usr/lib/ircd-hybrid/modules/autoload/m_ping.so
usr/lib/ircd-hybrid/modules/autoload/m_pong.la
usr/lib/ircd-hybrid/modules/autoload/m_pong.so
usr/lib/ircd-hybrid/modules/autoload/m_post.la
usr/lib/ircd-hybrid/modules/autoload/m_post.so
usr/lib/ircd-hybrid/modules/autoload/m_rehash.la
usr/lib/ircd-hybrid/modules/autoload/m_rehash.so
usr/lib/ircd-hybrid/modules/autoload/m_restart.la
usr/lib/ircd-hybrid/modules/autoload/m_restart.so
usr/lib/ircd-hybrid/modules/autoload/m_resv.la
usr/lib/ircd-hybrid/modules/autoload/m_resv.so
usr/lib/ircd-hybrid/modules/autoload/m_set.la
usr/lib/ircd-hybrid/modules/autoload/m_set.so
usr/lib/ircd-hybrid/modules/autoload/m_stats.la
usr/lib/ircd-hybrid/modules/autoload/m_stats.so
usr/lib/ircd-hybrid/modules/autoload/m_svinfo.la
usr/lib/ircd-hybrid/modules/autoload/m_svinfo.so
usr/lib/ircd-hybrid/modules/autoload/m_svsaccount.la
usr/lib/ircd-hybrid/modules/autoload/m_svsaccount.so
usr/lib/ircd-hybrid/modules/autoload/m_svshost.la
usr/lib/ircd-hybrid/modules/autoload/m_svshost.so
usr/lib/ircd-hybrid/modules/autoload/m_svsjoin.la
usr/lib/ircd-hybrid/modules/autoload/m_svsjoin.so
usr/lib/ircd-hybrid/modules/autoload/m_svskill.la
usr/lib/ircd-hybrid/modules/autoload/m_svskill.so
usr/lib/ircd-hybrid/modules/autoload/m_svsmode.la
usr/lib/ircd-hybrid/modules/autoload/m_svsmode.so
usr/lib/ircd-hybrid/modules/autoload/m_svsnick.la
usr/lib/ircd-hybrid/modules/autoload/m_svsnick.so
usr/lib/ircd-hybrid/modules/autoload/m_svspart.la
usr/lib/ircd-hybrid/modules/autoload/m_svspart.so
usr/lib/ircd-hybrid/modules/autoload/m_svstag.la
usr/lib/ircd-hybrid/modules/autoload/m_svstag.so
usr/lib/ircd-hybrid/modules/autoload/m_tburst.la
usr/lib/ircd-hybrid/modules/autoload/m_tburst.so
usr/lib/ircd-hybrid/modules/autoload/m_time.la
usr/lib/ircd-hybrid/modules/autoload/m_time.so
usr/lib/ircd-hybrid/modules/autoload/m_topic.la
usr/lib/ircd-hybrid/modules/autoload/m_topic.so
usr/lib/ircd-hybrid/modules/autoload/m_trace.la
usr/lib/ircd-hybrid/modules/autoload/m_trace.so
usr/lib/ircd-hybrid/modules/autoload/m_undline.la
usr/lib/ircd-hybrid/modules/autoload/m_undline.so
usr/lib/ircd-hybrid/modules/autoload/m_unkline.la
usr/lib/ircd-hybrid/modules/autoload/m_unkline.so
usr/lib/ircd-hybrid/modules/autoload/m_unresv.la
usr/lib/ircd-hybrid/modules/autoload/m_unresv.so
usr/lib/ircd-hybrid/modules/autoload/m_unxline.la
usr/lib/ircd-hybrid/modules/autoload/m_unxline.so
usr/lib/ircd-hybrid/modules/autoload/m_user.la
usr/lib/ircd-hybrid/modules/autoload/m_user.so
usr/lib/ircd-hybrid/modules/autoload/m_userhost.la
usr/lib/ircd-hybrid/modules/autoload/m_userhost.so
usr/lib/ircd-hybrid/modules/autoload/m_version.la
usr/lib/ircd-hybrid/modules/autoload/m_version.so
usr/lib/ircd-hybrid/modules/autoload/m_wallops.la
usr/lib/ircd-hybrid/modules/autoload/m_wallops.so
usr/lib/ircd-hybrid/modules/autoload/m_webirc.la
usr/lib/ircd-hybrid/modules/autoload/m_webirc.so
usr/lib/ircd-hybrid/modules/autoload/m_who.la
usr/lib/ircd-hybrid/modules/autoload/m_who.so
usr/lib/ircd-hybrid/modules/autoload/m_whois.la
usr/lib/ircd-hybrid/modules/autoload/m_whois.so
usr/lib/ircd-hybrid/modules/autoload/m_whowas.la
usr/lib/ircd-hybrid/modules/autoload/m_whowas.so
usr/lib/ircd-hybrid/modules/autoload/m_xline.la
usr/lib/ircd-hybrid/modules/autoload/m_xline.so
usr/lib/ircd-hybrid/modules/extra/
usr/lib/ircd-hybrid/modules/extra/m_opme.la
usr/lib/ircd-hybrid/modules/extra/m_opme.so
usr/lib/ircd-hybrid/modules/m_bmask.la
usr/lib/ircd-hybrid/modules/m_bmask.so
usr/lib/ircd-hybrid/modules/m_die.la
usr/lib/ircd-hybrid/modules/m_die.so
usr/lib/ircd-hybrid/modules/m_error.la
usr/lib/ircd-hybrid/modules/m_error.so
usr/lib/ircd-hybrid/modules/m_join.la
usr/lib/ircd-hybrid/modules/m_join.so
usr/lib/ircd-hybrid/modules/m_kick.la
usr/lib/ircd-hybrid/modules/m_kick.so
usr/lib/ircd-hybrid/modules/m_kill.la
usr/lib/ircd-hybrid/modules/m_kill.so
usr/lib/ircd-hybrid/modules/m_message.la
usr/lib/ircd-hybrid/modules/m_message.so
usr/lib/ircd-hybrid/modules/m_mode.la
usr/lib/ircd-hybrid/modules/m_mode.so
usr/lib/ircd-hybrid/modules/m_nick.la
usr/lib/ircd-hybrid/modules/m_nick.so
usr/lib/ircd-hybrid/modules/m_part.la
usr/lib/ircd-hybrid/modules/m_part.so
usr/lib/ircd-hybrid/modules/m_quit.la
usr/lib/ircd-hybrid/modules/m_quit.so
usr/lib/ircd-hybrid/modules/m_server.la
usr/lib/ircd-hybrid/modules/m_server.so
usr/lib/ircd-hybrid/modules/m_sjoin.la
usr/lib/ircd-hybrid/modules/m_sjoin.so
usr/lib/ircd-hybrid/modules/m_squit.la
usr/lib/ircd-hybrid/modules/m_squit.so
usr/lib/ircd-hybrid/modules/m_tmode.la
usr/lib/ircd-hybrid/modules/m_tmode.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ircd-hybrid.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/ircd-hybrid.conf
usr/share/
usr/share/ircd-hybrid/
usr/share/ircd-hybrid/help/
usr/share/ircd-hybrid/help/accept
usr/share/ircd-hybrid/help/admin
usr/share/ircd-hybrid/help/away
usr/share/ircd-hybrid/help/close
usr/share/ircd-hybrid/help/cmode
usr/share/ircd-hybrid/help/connect
usr/share/ircd-hybrid/help/die
usr/share/ircd-hybrid/help/dline
usr/share/ircd-hybrid/help/etrace
usr/share/ircd-hybrid/help/extban
usr/share/ircd-hybrid/help/hash
usr/share/ircd-hybrid/help/help
usr/share/ircd-hybrid/help/index
usr/share/ircd-hybrid/help/info
usr/share/ircd-hybrid/help/invite
usr/share/ircd-hybrid/help/ison
usr/share/ircd-hybrid/help/join
usr/share/ircd-hybrid/help/kick
usr/share/ircd-hybrid/help/kill
usr/share/ircd-hybrid/help/kline
usr/share/ircd-hybrid/help/knock
usr/share/ircd-hybrid/help/links
usr/share/ircd-hybrid/help/list
usr/share/ircd-hybrid/help/locops
usr/share/ircd-hybrid/help/lusers
usr/share/ircd-hybrid/help/map
usr/share/ircd-hybrid/help/module
usr/share/ircd-hybrid/help/monitor
usr/share/ircd-hybrid/help/motd
usr/share/ircd-hybrid/help/names
usr/share/ircd-hybrid/help/nick
usr/share/ircd-hybrid/help/notice
usr/share/ircd-hybrid/help/oper
usr/share/ircd-hybrid/help/part
usr/share/ircd-hybrid/help/pass
usr/share/ircd-hybrid/help/ping
usr/share/ircd-hybrid/help/pong
usr/share/ircd-hybrid/help/post
usr/share/ircd-hybrid/help/privmsg
usr/share/ircd-hybrid/help/quit
usr/share/ircd-hybrid/help/rehash
usr/share/ircd-hybrid/help/restart
usr/share/ircd-hybrid/help/resv
usr/share/ircd-hybrid/help/set
usr/share/ircd-hybrid/help/squit
usr/share/ircd-hybrid/help/stats
usr/share/ircd-hybrid/help/time
usr/share/ircd-hybrid/help/topic
usr/share/ircd-hybrid/help/trace
usr/share/ircd-hybrid/help/umode
usr/share/ircd-hybrid/help/undline
usr/share/ircd-hybrid/help/unkline
usr/share/ircd-hybrid/help/unresv
usr/share/ircd-hybrid/help/unxline
usr/share/ircd-hybrid/help/user
usr/share/ircd-hybrid/help/userhost
usr/share/ircd-hybrid/help/version
usr/share/ircd-hybrid/help/wallops
usr/share/ircd-hybrid/help/who
usr/share/ircd-hybrid/help/whois
usr/share/ircd-hybrid/help/whowas
usr/share/ircd-hybrid/help/xline
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/ircd.8.gz
var/
var/lib/
var/lib/ircd/
var/log/
var/log/ircd/
var/run/
