# files/certmonger-0.79.14-1/files
%FILES%
etc/
etc/certmonger/
etc/certmonger/certmonger.conf
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/certmonger.conf
usr/
usr/bin/
usr/bin/certmonger
usr/bin/getcert
usr/bin/ipa-getcert
usr/bin/local-getcert
usr/bin/selfsign-getcert
usr/lib/
usr/lib/certmonger/
usr/lib/certmonger/certmonger/
usr/lib/certmonger/certmonger/certmonger-session
usr/lib/certmonger/certmonger/dogtag-ipa-renew-agent-submit
usr/lib/certmonger/certmonger/dogtag-submit
usr/lib/certmonger/certmonger/ipa-submit
usr/lib/certmonger/certmonger/local-submit
usr/lib/certmonger/certmonger/scep-submit
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/certmonger.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/certmonger.conf
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/services/
usr/share/dbus-1/services/certmonger.service
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.fedorahosted.certmonger.service
usr/share/locale/
usr/share/locale/bg/
usr/share/locale/bg/LC_MESSAGES/
usr/share/locale/bg/LC_MESSAGES/certmonger.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/certmonger.mo
usr/share/locale/cs_CZ/
usr/share/locale/cs_CZ/LC_MESSAGES/
usr/share/locale/cs_CZ/LC_MESSAGES/certmonger.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/certmonger.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/certmonger.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/certmonger.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/certmonger.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/certmonger.mo
usr/share/locale/gu/
usr/share/locale/gu/LC_MESSAGES/
usr/share/locale/gu/LC_MESSAGES/certmonger.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/certmonger.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/certmonger.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/certmonger.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/certmonger.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/certmonger.mo
usr/share/locale/lt/
usr/share/locale/lt/LC_MESSAGES/
usr/share/locale/lt/LC_MESSAGES/certmonger.mo
usr/share/locale/nb/
usr/share/locale/nb/LC_MESSAGES/
usr/share/locale/nb/LC_MESSAGES/certmonger.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/certmonger.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/certmonger.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/certmonger.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/certmonger.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/certmonger.mo
usr/share/locale/si/
usr/share/locale/si/LC_MESSAGES/
usr/share/locale/si/LC_MESSAGES/certmonger.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/certmonger.mo
usr/share/locale/ta/
usr/share/locale/ta/LC_MESSAGES/
usr/share/locale/ta/LC_MESSAGES/certmonger.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/certmonger.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/certmonger.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/certmonger.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/certmonger.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/getcert-add-ca.1.gz
usr/share/man/man1/getcert-add-scep-ca.1.gz
usr/share/man/man1/getcert-list-cas.1.gz
usr/share/man/man1/getcert-list.1.gz
usr/share/man/man1/getcert-modify-ca.1.gz
usr/share/man/man1/getcert-refresh-ca.1.gz
usr/share/man/man1/getcert-refresh.1.gz
usr/share/man/man1/getcert-rekey.1.gz
usr/share/man/man1/getcert-remove-ca.1.gz
usr/share/man/man1/getcert-request.1.gz
usr/share/man/man1/getcert-resubmit.1.gz
usr/share/man/man1/getcert-start-tracking.1.gz
usr/share/man/man1/getcert-status.1.gz
usr/share/man/man1/getcert-stop-tracking.1.gz
usr/share/man/man1/getcert.1.gz
usr/share/man/man1/ipa-getcert.1.gz
usr/share/man/man1/local-getcert.1.gz
usr/share/man/man1/selfsign-getcert.1.gz
usr/share/man/man5/
usr/share/man/man5/certmonger.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/certmonger-dogtag-ipa-renew-agent-submit.8.gz
usr/share/man/man8/certmonger-dogtag-submit.8.gz
usr/share/man/man8/certmonger-ipa-submit.8.gz
usr/share/man/man8/certmonger-local-submit.8.gz
usr/share/man/man8/certmonger-scep-submit.8.gz
usr/share/man/man8/certmonger.8.gz
var/
var/lib/
var/lib/certmonger/
var/lib/certmonger/cas/
var/lib/certmonger/local/
var/lib/certmonger/requests/
