# files/btrfs-snapshots-git-0.6.1.r43.2d353ff-1/files
%FILES%
usr/
usr/bin/
usr/bin/btrfs-snapshots
usr/lib/
usr/lib/btrfs-snapshots/
usr/lib/btrfs-snapshots/action.bash
usr/lib/btrfs-snapshots/config.bash
usr/lib/btrfs-snapshots/event.bash
usr/lib/btrfs-snapshots/init.bash
usr/lib/btrfs-snapshots/locale.bash
usr/lib/btrfs-snapshots/snapshot.bash
usr/lib/btrfs-snapshots/timestamp.bash
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/btrfs-snapshots-create.service
usr/lib/systemd/system/btrfs-snapshots-create.timer
usr/lib/systemd/system/btrfs-snapshots-prune.service
usr/lib/systemd/system/btrfs-snapshots-prune.timer
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/btrfs-snapshots
usr/share/doc/
usr/share/doc/btrfs-snapshots/
usr/share/doc/btrfs-snapshots/LICENSE.md
usr/share/doc/btrfs-snapshots/README.md
usr/share/doc/btrfs-snapshots/config/
usr/share/doc/btrfs-snapshots/config/defaults.conf
usr/share/doc/btrfs-snapshots/config/profile.conf
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/btrfs-snapshots.5.gz
usr/share/man/man8/
usr/share/man/man8/btrfs-snapshots.8.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_btrfs-snapshots
