# files/linux-wifi-hotspot-4.0.3-1/files
%FILES%
etc/
etc/create_ap.conf
usr/
usr/bin/
usr/bin/create_ap
usr/bin/wihotspot
usr/bin/wihotspot-gui
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/create_ap.service
usr/share/
usr/share/applications/
usr/share/applications/wihotspot.desktop
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/create_ap
usr/share/doc/
usr/share/doc/create_ap/
usr/share/doc/create_ap/README.md
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/256x256/
usr/share/icons/hicolor/256x256/apps/
usr/share/icons/hicolor/256x256/apps/wihotspot.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/wihotspot.png
usr/share/icons/hicolor/512x512/
usr/share/icons/hicolor/512x512/apps/
usr/share/icons/hicolor/512x512/apps/wihotspot.png
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/apps/
usr/share/icons/hicolor/64x64/apps/wihotspot.png
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/wihotspot.svg
usr/share/licenses/
usr/share/licenses/linux-wifi-hotspot/
usr/share/licenses/linux-wifi-hotspot/LICENSE
usr/share/pixmaps/
usr/share/pixmaps/wihotspot.png
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/org.opensuse.policykit.wihotspot.policy
usr/share/polkit-1/rules.d/
usr/share/polkit-1/rules.d/90-org.opensuse.policykit.wihotspot.rules
