# files/hblock-3.3.1-1/files
%FILES%
usr/
usr/bin/
usr/bin/hblock
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hblock.service
usr/lib/systemd/system/hblock.timer
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/hblock.1.gz
