# files/nats-server-2.9.2-1/files
%FILES%
usr/
usr/bin/
usr/bin/nats-server
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nats-server.service
usr/share/
usr/share/doc/
usr/share/doc/nats-server/
usr/share/doc/nats-server/README.md
