# files/find-the-command-git-1.4.2.r2.g0fc7e85-2/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pacman-files.service
usr/lib/systemd/system/pacman-files.timer
usr/share/
usr/share/doc/
usr/share/doc/find-the-command/
usr/share/doc/find-the-command/README.md
usr/share/doc/find-the-command/ftc.bash
usr/share/doc/find-the-command/ftc.fish
usr/share/doc/find-the-command/ftc.zsh
usr/share/licenses/
usr/share/licenses/find-the-command/
usr/share/licenses/find-the-command/LICENSE
