# files/replay-sorcery-0.6.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/replay-sorcery
usr/etc/
usr/etc/replay-sorcery.conf
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/replay-sorcery-kms.service
usr/lib/systemd/user/
usr/lib/systemd/user/replay-sorcery.service
