# files/drkonqi-git-5.24.80_r777.g1fffaf6d-1/files
%FILES%
usr/
usr/bin/
usr/bin/drkonqi-coredump-gui
usr/lib/
usr/lib/drkonqi
usr/lib/drkonqi-coredump-cleanup
usr/lib/drkonqi-coredump-launcher
usr/lib/drkonqi-coredump-processor
usr/lib/qt/
usr/lib/qt/plugins/
usr/lib/qt/plugins/drkonqi/
usr/lib/qt/plugins/drkonqi/KDECoredumpNotifierTruck.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/drkonqi-coredump-processor@.service
usr/lib/systemd/user/
usr/lib/systemd/user/drkonqi-coredump-cleanup.service
usr/lib/systemd/user/drkonqi-coredump-cleanup.timer
usr/lib/systemd/user/drkonqi-coredump-launcher.socket
usr/lib/systemd/user/drkonqi-coredump-launcher@.service
usr/share/
usr/share/applications/
usr/share/applications/org.kde.drkonqi.coredump.gui.desktop
usr/share/applications/org.kde.drkonqi.desktop
usr/share/drkonqi/
usr/share/drkonqi/debuggers/
usr/share/drkonqi/debuggers/external.mac/
usr/share/drkonqi/debuggers/external.mac/lldbrc
usr/share/drkonqi/debuggers/external/
usr/share/drkonqi/debuggers/external/cdbrc
usr/share/drkonqi/debuggers/external/gdbrc
usr/share/drkonqi/debuggers/external/kdbgrc
usr/share/drkonqi/debuggers/external/lldbrc
usr/share/drkonqi/debuggers/internal/
usr/share/drkonqi/debuggers/internal/cdbrc
usr/share/drkonqi/debuggers/internal/dbxrc
usr/share/drkonqi/debuggers/internal/gdbrc
usr/share/drkonqi/debuggers/internal/kdbgwinrc
usr/share/drkonqi/debuggers/internal/lldbrc
usr/share/drkonqi/gdb/
usr/share/drkonqi/gdb/preamble.py
usr/share/drkonqi/mappings
usr/share/qlogging-categories5/
usr/share/qlogging-categories5/drkonqi.categories
