# files/systemd-cron-next-git-v1.0.2.r16.ga4c2de7-1/files
%FILES%
usr/
usr/bin/
usr/bin/crontab
usr/lib/
usr/lib/systemd-cron/
usr/lib/systemd-cron/boot-delay
usr/lib/systemd-cron/mail-on-failure
usr/lib/systemd-cron/remove-stale-stamps
usr/lib/systemd/
usr/lib/systemd/system-generators/
usr/lib/systemd/system-generators/systemd-crontab-generator
usr/lib/systemd/system/
usr/lib/systemd/system/cron-boot.service
usr/lib/systemd/system/cron-boot.target
usr/lib/systemd/system/cron-boot.timer
usr/lib/systemd/system/cron-daily.service
usr/lib/systemd/system/cron-daily.target
usr/lib/systemd/system/cron-daily.timer
usr/lib/systemd/system/cron-failure@.service
usr/lib/systemd/system/cron-hourly.service
usr/lib/systemd/system/cron-hourly.target
usr/lib/systemd/system/cron-hourly.timer
usr/lib/systemd/system/cron-monthly.service
usr/lib/systemd/system/cron-monthly.target
usr/lib/systemd/system/cron-monthly.timer
usr/lib/systemd/system/cron-update.path
usr/lib/systemd/system/cron-update.service
usr/lib/systemd/system/cron-weekly.service
usr/lib/systemd/system/cron-weekly.target
usr/lib/systemd/system/cron-weekly.timer
usr/lib/systemd/system/cron.target
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/crontab.1.gz
usr/share/man/man5/
usr/share/man/man5/anacrontab.5.gz
usr/share/man/man5/crontab.5.gz
usr/share/man/man7/
usr/share/man/man7/systemd.cron.7.gz
usr/share/man/man8/
usr/share/man/man8/systemd-crontab-generator.8.gz
var/
var/spool/
var/spool/cron/
