# files/nohang-0.2.0-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/nohang
etc/nohang/
etc/nohang/nohang-desktop.conf
etc/nohang/nohang.conf
usr/
usr/bin/
usr/bin/nohang
usr/bin/oom-sort
usr/bin/psi-top
usr/bin/psi2log
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nohang-desktop.service
usr/lib/systemd/system/nohang.service
usr/share/
usr/share/doc/
usr/share/doc/nohang/
usr/share/doc/nohang/CHANGELOG.md
usr/share/doc/nohang/README.md
usr/share/licenses/
usr/share/licenses/nohang/
usr/share/licenses/nohang/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/oom-sort.1.gz
usr/share/man/man1/psi-top.1.gz
usr/share/man/man1/psi2log.1.gz
usr/share/man/man8/
usr/share/man/man8/nohang.8.gz
usr/share/nohang/
usr/share/nohang/nohang-desktop.conf
usr/share/nohang/nohang.conf
usr/share/nohang/version
