# files/powerpanel-1.4.1-1/files
%FILES%
etc/
etc/hibernate.sh
etc/powerpanel/
etc/powerpanel/hibernate.sh
etc/powerpanel/pwrstatd-email.sh
etc/powerpanel/pwrstatd-lowbatt.sh
etc/powerpanel/pwrstatd-powerfail.sh
etc/powerpanel/shutdown.sh
etc/pwrstatd.conf
etc/shutdown.sh
usr/
usr/bin/
usr/bin/pwrstat
usr/bin/pwrstatd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pwrstatd.service
usr/share/
usr/share/licenses/
usr/share/licenses/powerpanel/
usr/share/licenses/powerpanel/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/pwrstat.8.gz
usr/share/man/man8/pwrstatd.8.gz
