# files/garuda-hotfixes-1.0.1-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/garuda-system-maintenance-wireplumber.service
usr/lib/systemd/system/garuda-system-maintenance-wireplumber.timer
