# files/neard-0.16-2/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/org.neard.conf
etc/neard/
etc/neard/main.conf
usr/
usr/bin/
usr/bin/nciattach
usr/bin/neard
usr/bin/nfctool
usr/include/
usr/include/near/
usr/include/near/adapter.h
usr/include/near/device.h
usr/include/near/log.h
usr/include/near/ndef.h
usr/include/near/nfc_copy.h
usr/include/near/plugin.h
usr/include/near/setting.h
usr/include/near/snep.h
usr/include/near/tag.h
usr/include/near/tlv.h
usr/include/near/types.h
usr/lib/
usr/lib/pkgconfig/
usr/lib/pkgconfig/neard.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/neard.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/nfctool.1.gz
usr/share/man/man5/
usr/share/man/man5/neard.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/neard.8.gz
