# files/tlp-git-1.5.0.r8.fce3f44-1/files
%FILES%
etc/
etc/tlp.conf
etc/tlp.d/
etc/tlp.d/00-template.conf
etc/tlp.d/README
usr/
usr/bin/
usr/bin/bluetooth
usr/bin/nfc
usr/bin/run-on-ac
usr/bin/run-on-bat
usr/bin/tlp
usr/bin/tlp-stat
usr/bin/wifi
usr/bin/wwan
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system-sleep/
usr/lib/systemd/system-sleep/tlp
usr/lib/systemd/system/
usr/lib/systemd/system/tlp.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/85-tlp.rules
usr/lib/udev/tlp-usb-udev
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/bluetooth
usr/share/bash-completion/completions/nfc
usr/share/bash-completion/completions/tlp
usr/share/bash-completion/completions/tlp-stat
usr/share/bash-completion/completions/wifi
usr/share/bash-completion/completions/wwan
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/bluetooth.1.gz
usr/share/man/man1/nfc.1.gz
usr/share/man/man1/run-on-ac.1.gz
usr/share/man/man1/run-on-bat.1.gz
usr/share/man/man1/wifi.1.gz
usr/share/man/man1/wwan.1.gz
usr/share/man/man8/
usr/share/man/man8/tlp-stat.8.gz
usr/share/man/man8/tlp.8.gz
usr/share/man/man8/tlp.service.8.gz
usr/share/metainfo/
usr/share/metainfo/de.linrunner.tlp.metainfo.xml
usr/share/tlp/
usr/share/tlp/bat.d/
usr/share/tlp/bat.d/05-thinkpad
usr/share/tlp/bat.d/10-thinkpad-legacy
usr/share/tlp/bat.d/15-asus
usr/share/tlp/bat.d/20-huawei
usr/share/tlp/bat.d/25-lenovo
usr/share/tlp/bat.d/30-samsung
usr/share/tlp/bat.d/35-lg
usr/share/tlp/bat.d/40-sony
usr/share/tlp/bat.d/90-generic
usr/share/tlp/bat.d/TEMPLATE
usr/share/tlp/defaults.conf
usr/share/tlp/func.d/
usr/share/tlp/func.d/05-tlp-func-pm
usr/share/tlp/func.d/10-tlp-func-cpu
usr/share/tlp/func.d/15-tlp-func-disk
usr/share/tlp/func.d/20-tlp-func-usb
usr/share/tlp/func.d/25-tlp-func-rf
usr/share/tlp/func.d/30-tlp-func-rf-sw
usr/share/tlp/func.d/35-tlp-func-batt
usr/share/tlp/func.d/40-tlp-func-bay
usr/share/tlp/func.d/45-tlp-func-gpu
usr/share/tlp/func.d/tlp-func-stat
usr/share/tlp/rename.conf
usr/share/tlp/tlp-func-base
usr/share/tlp/tlp-pcilist
usr/share/tlp/tlp-readconfs
usr/share/tlp/tlp-usblist
usr/share/tlp/tpacpi-bat
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_tlp
usr/share/zsh/site-functions/_tlp-radio-device
usr/share/zsh/site-functions/_tlp-run-on
usr/share/zsh/site-functions/_tlp-stat
var/
var/lib/
var/lib/tlp/
