# files/zram-generator-git-1.1.2+2+g1f82211-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system-generators/
usr/lib/systemd/system-generators/zram-generator
usr/lib/systemd/system/
usr/lib/systemd/system/systemd-zram-setup@.service
usr/share/
usr/share/doc/
usr/share/doc/zram-generator/
usr/share/doc/zram-generator/half-memory.conf.example
usr/share/doc/zram-generator/zram-generator.conf.example
usr/share/licenses/
usr/share/licenses/zram-generator-git/
usr/share/licenses/zram-generator-git/LICENSE
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/zram-generator.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/zram-generator.8.gz
