# files/rustdesk-server-bin-1.1.6-2/files
%FILES%
opt/
opt/rustdesk-server/
opt/rustdesk-server/data/
opt/rustdesk-server/hbbr
opt/rustdesk-server/hbbs
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rustdesk-server-hbbr.service
usr/lib/systemd/system/rustdesk-server-hbbs.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/rustdesk-server.conf
