# files/acestream-engine-3.1.49-3/files
%FILES%
opt/
opt/acestream/
opt/acestream/acestream.conf
opt/acestream/acestreamengine
opt/acestream/data/
opt/acestream/data/images/
opt/acestream/data/images/acestream16.png
opt/acestream/data/images/acestream32.png
opt/acestream/data/images/icon.ico
opt/acestream/data/images/streamer-16.png
opt/acestream/data/images/streamer-32.png
opt/acestream/data/images/user_pro.png
opt/acestream/data/lang/
opt/acestream/data/lang/en_EN.lang
opt/acestream/data/lang/ru_RU.lang
opt/acestream/data/lang/vi_VI.lang
opt/acestream/data/plugins/
opt/acestream/data/plugins/adultswim.py
opt/acestream/data/plugins/afreeca.py
opt/acestream/data/plugins/afreecatv.py
opt/acestream/data/plugins/aftonbladet.py
opt/acestream/data/plugins/alieztv.py
opt/acestream/data/plugins/antenna.py
opt/acestream/data/plugins/ard_live.py
opt/acestream/data/plugins/ard_mediathek.py
opt/acestream/data/plugins/artetv.py
opt/acestream/data/plugins/atresplayer.py
opt/acestream/data/plugins/azubutv.py
opt/acestream/data/plugins/bambuser.py
opt/acestream/data/plugins/beam.py
opt/acestream/data/plugins/beattv.py
opt/acestream/data/plugins/bigo.py
opt/acestream/data/plugins/bilibili.py
opt/acestream/data/plugins/bliptv.py
opt/acestream/data/plugins/chaturbate.py
opt/acestream/data/plugins/cinergroup.py
opt/acestream/data/plugins/common_jwplayer.py
opt/acestream/data/plugins/common_swf.py
opt/acestream/data/plugins/connectcast.py
opt/acestream/data/plugins/crunchyroll.py
opt/acestream/data/plugins/cybergame.py
opt/acestream/data/plugins/dailymotion.py
opt/acestream/data/plugins/dingittv.py
opt/acestream/data/plugins/disney_de.py
opt/acestream/data/plugins/dmcloud.py
opt/acestream/data/plugins/dmcloud_embed.py
opt/acestream/data/plugins/dogan.py
opt/acestream/data/plugins/dogus.py
opt/acestream/data/plugins/dommune.py
opt/acestream/data/plugins/douyutv.py
opt/acestream/data/plugins/dplay.py
opt/acestream/data/plugins/drdk.py
opt/acestream/data/plugins/euronews.py
opt/acestream/data/plugins/expressen.py
opt/acestream/data/plugins/facebook.py
opt/acestream/data/plugins/filmon.py
opt/acestream/data/plugins/filmon_us.py
opt/acestream/data/plugins/foxtr.py
opt/acestream/data/plugins/furstream.py
opt/acestream/data/plugins/gaminglive.py
opt/acestream/data/plugins/gomexp.py
opt/acestream/data/plugins/goodgame.py
opt/acestream/data/plugins/hitbox.py
opt/acestream/data/plugins/ine.py
opt/acestream/data/plugins/itvplayer.py
opt/acestream/data/plugins/kanal7.py
opt/acestream/data/plugins/letontv.py
opt/acestream/data/plugins/life_ru.py
opt/acestream/data/plugins/livecodingtv.py
opt/acestream/data/plugins/liveme.py
opt/acestream/data/plugins/livestation.py
opt/acestream/data/plugins/livestream.py
opt/acestream/data/plugins/media_ccc_de.py
opt/acestream/data/plugins/mediaklikk.py
opt/acestream/data/plugins/meerkat.py
opt/acestream/data/plugins/mips.py
opt/acestream/data/plugins/mlgtv.py
opt/acestream/data/plugins/news.py
opt/acestream/data/plugins/nhkworld.py
opt/acestream/data/plugins/nineanime.py
opt/acestream/data/plugins/nos.py
opt/acestream/data/plugins/npo.py
opt/acestream/data/plugins/nrk.py
opt/acestream/data/plugins/oldlivestream.py
opt/acestream/data/plugins/openrectv.py
opt/acestream/data/plugins/orf_tvthek.py
opt/acestream/data/plugins/ovva.py
opt/acestream/data/plugins/pandatv.py
opt/acestream/data/plugins/periscope.py
opt/acestream/data/plugins/picarto.py
opt/acestream/data/plugins/piczel.py
opt/acestream/data/plugins/powerapp.py
opt/acestream/data/plugins/rtlxl.py
opt/acestream/data/plugins/rtve.py
opt/acestream/data/plugins/ruv.py
opt/acestream/data/plugins/sbsdiscovery.py
opt/acestream/data/plugins/seemeplay.py
opt/acestream/data/plugins/servustv.py
opt/acestream/data/plugins/shakhtar.py
opt/acestream/data/plugins/speedrunslive.py
opt/acestream/data/plugins/sportschau.py
opt/acestream/data/plugins/ssh101.py
opt/acestream/data/plugins/stream.py
opt/acestream/data/plugins/streamboat.py
opt/acestream/data/plugins/streamingvideoprovider.py
opt/acestream/data/plugins/streamlive.py
opt/acestream/data/plugins/streamme.py
opt/acestream/data/plugins/streamupcom.py
opt/acestream/data/plugins/svtplay.py
opt/acestream/data/plugins/tga.py
opt/acestream/data/plugins/tigerdile.py
opt/acestream/data/plugins/trt.py
opt/acestream/data/plugins/turkuvaz.py
opt/acestream/data/plugins/tv360.py
opt/acestream/data/plugins/tv3cat.py
opt/acestream/data/plugins/tv4play.py
opt/acestream/data/plugins/tv8.py
opt/acestream/data/plugins/tvcatchup.py
opt/acestream/data/plugins/tvplayer.py
opt/acestream/data/plugins/twitch.py
opt/acestream/data/plugins/ustreamtv.py
opt/acestream/data/plugins/vaughnlive.py
opt/acestream/data/plugins/veetle.py
opt/acestream/data/plugins/vgtv.py
opt/acestream/data/plugins/viagame.py
opt/acestream/data/plugins/viasat.py
opt/acestream/data/plugins/viasat_embed.py
opt/acestream/data/plugins/vidio.py
opt/acestream/data/plugins/vk.py
opt/acestream/data/plugins/wattv.py
opt/acestream/data/plugins/webtv.py
opt/acestream/data/plugins/weeb.py
opt/acestream/data/plugins/younow.py
opt/acestream/data/plugins/youtube.py
opt/acestream/data/plugins/zdf_mediathek.py
opt/acestream/data/schema_sdb_v6.sql
opt/acestream/data/webui/
opt/acestream/data/webui/cert/
opt/acestream/data/webui/cert/engine.pem
opt/acestream/data/webui/css/
opt/acestream/data/webui/css/ace-style.css
opt/acestream/data/webui/css/bootstrap-datetimepicker.min.css
opt/acestream/data/webui/css/bootstrap-responsive.css
opt/acestream/data/webui/css/bootstrap.css
opt/acestream/data/webui/css/common.css
opt/acestream/data/webui/css/css.css
opt/acestream/data/webui/css/font-awesome.min.css
opt/acestream/data/webui/css/gritter.css
opt/acestream/data/webui/css/ie.css
opt/acestream/data/webui/css/jquery-ui.css
opt/acestream/data/webui/css/jquery.mCustomScrollbar.css
opt/acestream/data/webui/css/light.css
opt/acestream/data/webui/css/media.css
opt/acestream/data/webui/css/popup.css
opt/acestream/data/webui/css/popup2.css
opt/acestream/data/webui/css/popup3.css
opt/acestream/data/webui/css/remote-control.css
opt/acestream/data/webui/css/scroll.css
opt/acestream/data/webui/css/seekbar.css
opt/acestream/data/webui/css/slide-menu.css
opt/acestream/data/webui/css/style.css
opt/acestream/data/webui/css/stylesheet.css
opt/acestream/data/webui/css/ui-menu.css
opt/acestream/data/webui/css/unicorn.css
opt/acestream/data/webui/css/uniform.css
opt/acestream/data/webui/fonts/
opt/acestream/data/webui/fonts/fontawesome-webfont.eot
opt/acestream/data/webui/fonts/fontawesome-webfont.svg
opt/acestream/data/webui/fonts/fontawesome-webfont.ttf
opt/acestream/data/webui/fonts/fontawesome-webfont.woff
opt/acestream/data/webui/fonts/fontawesome-webfont.woff2
opt/acestream/data/webui/html/
opt/acestream/data/webui/html/download.html
opt/acestream/data/webui/html/iframe_playlist_add_item.html
opt/acestream/data/webui/html/monitor.html
opt/acestream/data/webui/html/player.html
opt/acestream/data/webui/html/remote_control.html
opt/acestream/data/webui/html/server.html
opt/acestream/data/webui/html/server_auth.html
opt/acestream/data/webui/images/
opt/acestream/data/webui/images/arrows.png
opt/acestream/data/webui/images/breadcrumb.png
opt/acestream/data/webui/images/close.png
opt/acestream/data/webui/images/close_36.png
opt/acestream/data/webui/images/favicon.ico
opt/acestream/data/webui/images/glyphicons-halflings-white.png
opt/acestream/data/webui/images/glyphicons-halflings.png
opt/acestream/data/webui/images/headBG.png
opt/acestream/data/webui/images/icon-broadcast.png
opt/acestream/data/webui/images/icon-copy.png
opt/acestream/data/webui/images/icon-premium.png
opt/acestream/data/webui/images/icon-proxy.png
opt/acestream/data/webui/images/icon-update.png
opt/acestream/data/webui/images/icon-view.png
opt/acestream/data/webui/images/loader.gif
opt/acestream/data/webui/images/logo-small.png
opt/acestream/data/webui/images/logo.png
opt/acestream/data/webui/images/logo2.png
opt/acestream/data/webui/images/remote-control/
opt/acestream/data/webui/images/remote-control/background.png
opt/acestream/data/webui/images/remote-control/ic_pause_white_48dp.png
opt/acestream/data/webui/images/remote-control/ic_play_arrow_white_48dp.png
opt/acestream/data/webui/images/remote-control/rc_crop_disabled.png
opt/acestream/data/webui/images/remote-control/rc_crop_normal.png
opt/acestream/data/webui/images/remote-control/rc_crop_pressed.png
opt/acestream/data/webui/images/remote-control/rc_image_audio.png
opt/acestream/data/webui/images/remote-control/rc_image_blank.png
opt/acestream/data/webui/images/remote-control/rc_image_live.png
opt/acestream/data/webui/images/remote-control/rc_image_video.png
opt/acestream/data/webui/images/remote-control/rc_next_disabled.png
opt/acestream/data/webui/images/remote-control/rc_next_normal.png
opt/acestream/data/webui/images/remote-control/rc_next_pressed.png
opt/acestream/data/webui/images/remote-control/rc_pause_100.png
opt/acestream/data/webui/images/remote-control/rc_pause_25.png
opt/acestream/data/webui/images/remote-control/rc_pause_75.png
opt/acestream/data/webui/images/remote-control/rc_play_100.png
opt/acestream/data/webui/images/remote-control/rc_play_25.png
opt/acestream/data/webui/images/remote-control/rc_play_75.png
opt/acestream/data/webui/images/remote-control/rc_prev_disabled.png
opt/acestream/data/webui/images/remote-control/rc_prev_normal.png
opt/acestream/data/webui/images/remote-control/rc_prev_pressed.png
opt/acestream/data/webui/images/remote-control/rc_restart_disabled.png
opt/acestream/data/webui/images/remote-control/rc_restart_normal.png
opt/acestream/data/webui/images/remote-control/rc_restart_pressed.png
opt/acestream/data/webui/images/remote-control/rc_select_audio_track_disabled.png
opt/acestream/data/webui/images/remote-control/rc_select_audio_track_normal.png
opt/acestream/data/webui/images/remote-control/rc_select_audio_track_pressed.png
opt/acestream/data/webui/images/remote-control/rc_select_device_disabled.png
opt/acestream/data/webui/images/remote-control/rc_select_device_normal.png
opt/acestream/data/webui/images/remote-control/rc_select_device_pressed.png
opt/acestream/data/webui/images/remote-control/rc_select_subtitle_track_disabled.png
opt/acestream/data/webui/images/remote-control/rc_select_subtitle_track_normal.png
opt/acestream/data/webui/images/remote-control/rc_select_subtitle_track_pressed.png
opt/acestream/data/webui/images/remote-control/rc_stop_disabled.png
opt/acestream/data/webui/images/remote-control/rc_stop_normal.png
opt/acestream/data/webui/images/remote-control/rc_stop_pressed.png
opt/acestream/data/webui/images/remote-control/rc_volume_max.png
opt/acestream/data/webui/images/remote-control/rc_volume_min.png
opt/acestream/data/webui/images/remote-control/ring-92.gif
opt/acestream/data/webui/javascript/
opt/acestream/data/webui/javascript/bootstrap-datetimepicker.min.js
opt/acestream/data/webui/javascript/bootstrap.js
opt/acestream/data/webui/javascript/broadcasting.js
opt/acestream/data/webui/javascript/clipboard.min.js
opt/acestream/data/webui/javascript/init.js
opt/acestream/data/webui/javascript/jquery-1.10.2.min.js
opt/acestream/data/webui/javascript/jquery-ui.min.js
opt/acestream/data/webui/javascript/jquery.mCustomScrollbar.concat.min.js
opt/acestream/data/webui/javascript/jquery.sticky.js
opt/acestream/data/webui/javascript/lang.en_EN.js
opt/acestream/data/webui/javascript/lang.js
opt/acestream/data/webui/javascript/lang.ru_RU.js
opt/acestream/data/webui/javascript/languages.json
opt/acestream/data/webui/javascript/moment.min.js
opt/acestream/data/webui/javascript/p2p-graph.min.js
opt/acestream/data/webui/javascript/proxy_server.js
opt/acestream/data/webui/javascript/q.js
opt/acestream/data/webui/javascript/remote-control.js
opt/acestream/data/webui/javascript/seekbar.js
opt/acestream/data/webui/javascript/server.js
opt/acestream/data/webui/javascript/strftime.min.js
opt/acestream/data/webui/javascript/uniform.js
opt/acestream/lib/
opt/acestream/lib/GeoIP-1.3.2-py2.7-linux-x86_64.egg
opt/acestream/lib/M2Crypto-0.35.2-py2.7-linux-x86_64.egg
opt/acestream/lib/Routes-2.4.1-py2.7.egg
opt/acestream/lib/acestreamengine/
opt/acestream/lib/acestreamengine/Core.so
opt/acestream/lib/acestreamengine/CoreApp.so
opt/acestream/lib/acestreamengine/Transport.so
opt/acestream/lib/acestreamengine/__init__.py
opt/acestream/lib/acestreamengine/live.so
opt/acestream/lib/acestreamengine/node.so
opt/acestream/lib/acestreamengine/streamer.so
opt/acestream/lib/beautifulsoup4-4.5.3-py2.7.egg
opt/acestream/lib/bitarray-0.8.1-py2.7-linux-x86_64.egg
opt/acestream/lib/blist-1.3.4-py2.7-linux-x86_64.egg
opt/acestream/lib/dnspython-1.15.0-py2.7.egg
opt/acestream/lib/enum34-1.1.6-py2.7.egg
opt/acestream/lib/futures-3.0.5-py2.7.egg
opt/acestream/lib/html5lib-0.999999999-py2.7.egg
opt/acestream/lib/ipaddr-2.2.0-py2.7.egg
opt/acestream/lib/ipwhois-1.0.0-py2.7.egg
opt/acestream/lib/iso8601-0.1.10-py2.7.egg
opt/acestream/lib/lxml-3.7.2-py2.7-linux-x86_64.egg
opt/acestream/lib/miniupnpc-2.1-py2.7-linux-x86_64.egg
opt/acestream/lib/netifaces-0.10.5-py2.7-linux-x86_64.egg
opt/acestream/lib/opengraph-0.5-py2.7.egg
opt/acestream/lib/protobuf-3.0.0b2-py2.7.egg
opt/acestream/lib/psutil-1.2.1-py2.7-linux-x86_64.egg
opt/acestream/lib/pycompat27.so
opt/acestream/lib/repoze.lru-0.6-py2.7.egg
opt/acestream/lib/requests-2.12.5-py2.7.egg
opt/acestream/lib/singledispatch-3.4.0.3-py2.7.egg
opt/acestream/lib/six-1.10.0-py2.7.egg
opt/acestream/lib/typing-3.7.4-py2.7.egg
opt/acestream/lib/webencodings-0.5-py2.7.egg
opt/acestream/lib/websocket_client-0.40.0-py2.7.egg
opt/acestream/lib/xmltodict-0.9.2-py2.7.egg
opt/acestream/start-engine
usr/
usr/bin/
usr/bin/acestreamengine
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/acestream-engine.service
