# files/j9-git-r2812.3c634dad-1/files
%FILES%
usr/
usr/bin/
usr/bin/j9
usr/bin/j9hsd
usr/lib/
usr/lib/j9/
usr/lib/j9/addons/
usr/lib/j9/addons/config/
usr/lib/j9/addons/config/addins.txt
usr/lib/j9/addons/ide/
usr/lib/j9/addons/ide/jhs/
usr/lib/j9/addons/ide/jhs/app/
usr/lib/j9/addons/ide/jhs/app/app1.ijs
usr/lib/j9/addons/ide/jhs/app/app2.ijs
usr/lib/j9/addons/ide/jhs/app/app3.ijs
usr/lib/j9/addons/ide/jhs/app/app4.ijs
usr/lib/j9/addons/ide/jhs/app/app5.ijs
usr/lib/j9/addons/ide/jhs/app/app6.ijs
usr/lib/j9/addons/ide/jhs/cojhs/
usr/lib/j9/addons/ide/jhs/cojhs/jd3.ijs
usr/lib/j9/addons/ide/jhs/cojhs/jdebug.ijs
usr/lib/j9/addons/ide/jhs/cojhs/jhot.ijs
usr/lib/j9/addons/ide/jhs/cojhs/jtable.ijs
usr/lib/j9/addons/ide/jhs/cojhs/jwatch.ijs
usr/lib/j9/addons/ide/jhs/config/
usr/lib/j9/addons/ide/jhs/config/jhs.cfg
usr/lib/j9/addons/ide/jhs/core.ijs
usr/lib/j9/addons/ide/jhs/d3.ijs
usr/lib/j9/addons/ide/jhs/debug.ijs
usr/lib/j9/addons/ide/jhs/demo/
usr/lib/j9/addons/ide/jhs/demo/jdemo1.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo10.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo11.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo12.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo13.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo14.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo2.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo3.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo4.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo5.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo6.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo7.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo8.ijs
usr/lib/j9/addons/ide/jhs/demo/jdemo9.ijs
usr/lib/j9/addons/ide/jhs/favicon.ico
usr/lib/j9/addons/ide/jhs/gnuplot.ijs
usr/lib/j9/addons/ide/jhs/gpdemo.ijs
usr/lib/j9/addons/ide/jhs/history.txt
usr/lib/j9/addons/ide/jhs/installer.ijs
usr/lib/j9/addons/ide/jhs/jal.ijs
usr/lib/j9/addons/ide/jhs/jbinfo.ijs
usr/lib/j9/addons/ide/jhs/jfif.ijs
usr/lib/j9/addons/ide/jhs/jfile.ijs
usr/lib/j9/addons/ide/jhs/jfiles.ijs
usr/lib/j9/addons/ide/jhs/jfilesrc.ijs
usr/lib/j9/addons/ide/jhs/jgcp.ijs
usr/lib/j9/addons/ide/jhs/jijs.ijs
usr/lib/j9/addons/ide/jhs/jijx.ijs
usr/lib/j9/addons/ide/jhs/jj.ijs
usr/lib/j9/addons/ide/jhs/jjserver.ijs
usr/lib/j9/addons/ide/jhs/jlogin.ijs
usr/lib/j9/addons/ide/jhs/js/
usr/lib/j9/addons/ide/jhs/js/codemirror/
usr/lib/j9/addons/ide/jhs/js/codemirror/codemirror.4.2.css
usr/lib/j9/addons/ide/jhs/js/codemirror/codemirror.min.4.2.js
usr/lib/j9/addons/ide/jhs/js/codemirror/j/
usr/lib/j9/addons/ide/jhs/js/codemirror/j/j.4.2.js
usr/lib/j9/addons/ide/jhs/js/codemirror/j/jtheme.4.2.css
usr/lib/j9/addons/ide/jhs/js/codemirror/license.txt
usr/lib/j9/addons/ide/jhs/js/codemirror/util/
usr/lib/j9/addons/ide/jhs/js/codemirror/util/dialog.4.2.css
usr/lib/j9/addons/ide/jhs/js/codemirror/util/dialog.4.2.js
usr/lib/j9/addons/ide/jhs/js/codemirror/util/search.4.2.js
usr/lib/j9/addons/ide/jhs/js/codemirror/util/searchcursor.4.2.js
usr/lib/j9/addons/ide/jhs/js/d3/
usr/lib/j9/addons/ide/jhs/js/d3/d3.css
usr/lib/j9/addons/ide/jhs/js/d3/d3.js
usr/lib/j9/addons/ide/jhs/js/d3/d3.min.js
usr/lib/j9/addons/ide/jhs/js/d3/license
usr/lib/j9/addons/ide/jhs/js/handsontable/
usr/lib/j9/addons/ide/jhs/js/handsontable/handsontable.full.min.css
usr/lib/j9/addons/ide/jhs/js/handsontable/handsontable.full.min.js
usr/lib/j9/addons/ide/jhs/js/jquery/
usr/lib/j9/addons/ide/jhs/js/jquery/jquery-2.0.3.min.js
usr/lib/j9/addons/ide/jhs/js/jquery/jquery-ui.css
usr/lib/j9/addons/ide/jhs/js/jquery/jquery-ui.min.js
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/animated-overlay.gif
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-bg_flat_0_aaaaaa_40x100.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-bg_flat_75_ffffff_40x100.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-bg_glass_55_fbf9ee_1x400.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-bg_glass_65_ffffff_1x400.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-bg_glass_75_dadada_1x400.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-bg_glass_75_e6e6e6_1x400.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-bg_glass_95_fef1ec_1x400.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-bg_highlight-soft_75_cccccc_1x100.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-icons_222222_256x240.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-icons_2e83ff_256x240.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-icons_454545_256x240.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-icons_888888_256x240.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/images/ui-icons_cd0a0a_256x240.png
usr/lib/j9/addons/ide/jhs/js/jquery/smoothness/jquery-ui.custom.css
usr/lib/j9/addons/ide/jhs/js/jsoftware/
usr/lib/j9/addons/ide/jhs/js/jsoftware/d3_basic.js
usr/lib/j9/addons/ide/jhs/js/webgl/
usr/lib/j9/addons/ide/jhs/js/webgl/glmatrix-min.js
usr/lib/j9/addons/ide/jhs/jum.ijs
usr/lib/j9/addons/ide/jhs/manifest.ijs
usr/lib/j9/addons/ide/jhs/misc/
usr/lib/j9/addons/ide/jhs/misc/avx.ijs
usr/lib/j9/addons/ide/jhs/readme.txt
usr/lib/j9/addons/ide/jhs/sp.ijs
usr/lib/j9/addons/ide/jhs/spx/
usr/lib/j9/addons/ide/jhs/spx/plot.ijs
usr/lib/j9/addons/ide/jhs/spx/spx.ijs
usr/lib/j9/addons/ide/jhs/tool.ijs
usr/lib/j9/addons/ide/jhs/util.ijs
usr/lib/j9/addons/ide/jhs/utilh.ijs
usr/lib/j9/addons/ide/jhs/utiljs.ijs
usr/lib/j9/addons/ide/qt/
usr/lib/j9/addons/ide/qt/config/
usr/lib/j9/addons/ide/qt/config/case.cfg
usr/lib/j9/addons/ide/qt/config/dirmatch.cfg
usr/lib/j9/addons/ide/qt/config/launch.cfg
usr/lib/j9/addons/ide/qt/config/userkeys.cfg
usr/lib/j9/addons/ide/qt/debugs.ijs
usr/lib/j9/addons/ide/qt/gl2.ijs
usr/lib/j9/addons/ide/qt/help.ijs
usr/lib/j9/addons/ide/qt/images/
usr/lib/j9/addons/ide/qt/images/about.png
usr/lib/j9/addons/ide/qt/images/clear.png
usr/lib/j9/addons/ide/qt/images/copy.png
usr/lib/j9/addons/ide/qt/images/cut.png
usr/lib/j9/addons/ide/qt/images/cutback.png
usr/lib/j9/addons/ide/qt/images/dissect-current.png
usr/lib/j9/addons/ide/qt/images/dissect-toggle.png
usr/lib/j9/addons/ide/qt/images/go-bottom.png
usr/lib/j9/addons/ide/qt/images/go-down.png
usr/lib/j9/addons/ide/qt/images/go-first.png
usr/lib/j9/addons/ide/qt/images/go-home.png
usr/lib/j9/addons/ide/qt/images/go-last.png
usr/lib/j9/addons/ide/qt/images/go-next.png
usr/lib/j9/addons/ide/qt/images/go-previous.png
usr/lib/j9/addons/ide/qt/images/go-top.png
usr/lib/j9/addons/ide/qt/images/go-up.png
usr/lib/j9/addons/ide/qt/images/goback.png
usr/lib/j9/addons/ide/qt/images/goforward.png
usr/lib/j9/addons/ide/qt/images/help-contents.png
usr/lib/j9/addons/ide/qt/images/home.png
usr/lib/j9/addons/ide/qt/images/jump-to.png
usr/lib/j9/addons/ide/qt/images/media-playback-pause.png
usr/lib/j9/addons/ide/qt/images/media-seek-forward.png
usr/lib/j9/addons/ide/qt/images/new.png
usr/lib/j9/addons/ide/qt/images/open.png
usr/lib/j9/addons/ide/qt/images/paste.png
usr/lib/j9/addons/ide/qt/images/print.png
usr/lib/j9/addons/ide/qt/images/process-stop.png
usr/lib/j9/addons/ide/qt/images/redo.png
usr/lib/j9/addons/ide/qt/images/refresh.png
usr/lib/j9/addons/ide/qt/images/run.png
usr/lib/j9/addons/ide/qt/images/runcursor.png
usr/lib/j9/addons/ide/qt/images/spreadsheet.png
usr/lib/j9/addons/ide/qt/images/stack.png
usr/lib/j9/addons/ide/qt/images/stepinto.png
usr/lib/j9/addons/ide/qt/images/stepout.png
usr/lib/j9/addons/ide/qt/images/stepover.png
usr/lib/j9/addons/ide/qt/images/stopmanager.png
usr/lib/j9/addons/ide/qt/images/stopname.png
usr/lib/j9/addons/ide/qt/images/system-run.png
usr/lib/j9/addons/ide/qt/images/undo.png
usr/lib/j9/addons/ide/qt/images/watchmanager.png
usr/lib/j9/addons/ide/qt/images/zoom-in.png
usr/lib/j9/addons/ide/qt/keys.ijs
usr/lib/j9/addons/ide/qt/manifest.ijs
usr/lib/j9/addons/ide/qt/opengl.ijs
usr/lib/j9/addons/ide/qt/pacman.ijs
usr/lib/j9/addons/ide/qt/printer.ijs
usr/lib/j9/addons/ide/qt/qt.ijs
usr/lib/j9/addons/ide/qt/qtlib.ijs
usr/lib/j9/bin/
usr/lib/j9/bin/icons/
usr/lib/j9/bin/icons/app.png
usr/lib/j9/bin/icons/favicon.ico
usr/lib/j9/bin/icons/j9gray.png
usr/lib/j9/bin/icons/j9gray.svg
usr/lib/j9/bin/icons/j9green.png
usr/lib/j9/bin/icons/j9green.svg
usr/lib/j9/bin/icons/j9red.png
usr/lib/j9/bin/icons/j9red.svg
usr/lib/j9/bin/icons/jblue.icns
usr/lib/j9/bin/icons/jblue.ico
usr/lib/j9/bin/icons/jblue.png
usr/lib/j9/bin/icons/jgray.icns
usr/lib/j9/bin/icons/jgray.ico
usr/lib/j9/bin/icons/jgray.png
usr/lib/j9/bin/icons/jgray.svg
usr/lib/j9/bin/icons/jgreen.icns
usr/lib/j9/bin/icons/jgreen.ico
usr/lib/j9/bin/icons/jgreen.png
usr/lib/j9/bin/icons/jgreen.svg
usr/lib/j9/bin/icons/jred.icns
usr/lib/j9/bin/icons/jred.ico
usr/lib/j9/bin/icons/jred.png
usr/lib/j9/bin/icons/jred.svg
usr/lib/j9/bin/icons/jyellow.icns
usr/lib/j9/bin/icons/jyellow.ico
usr/lib/j9/bin/icons/jyellow.png
usr/lib/j9/bin/installer.txt
usr/lib/j9/bin/jbrk
usr/lib/j9/bin/jconsole
usr/lib/j9/bin/jhs
usr/lib/j9/bin/jqt
usr/lib/j9/bin/libjqt.so
usr/lib/j9/bin/libjqt.so.1
usr/lib/j9/bin/libjqt.so.1.9
usr/lib/j9/bin/libjqt.so.1.9.5
usr/lib/j9/bin/libtsdll.so
usr/lib/j9/bin/profile.ijs
usr/lib/j9/bin/profilex.ijs
usr/lib/j9/bin/profilex_template.ijs
usr/lib/j9/system/
usr/lib/j9/system/config/
usr/lib/j9/system/config/base.cfg
usr/lib/j9/system/config/case.cfg
usr/lib/j9/system/config/folders.cfg
usr/lib/j9/system/config/version.txt
usr/lib/j9/system/defs/
usr/lib/j9/system/defs/hostdefs_aix.ijs
usr/lib/j9/system/defs/hostdefs_android.ijs
usr/lib/j9/system/defs/hostdefs_android_64.ijs
usr/lib/j9/system/defs/hostdefs_darwin.ijs
usr/lib/j9/system/defs/hostdefs_darwin_64.ijs
usr/lib/j9/system/defs/hostdefs_linux.ijs
usr/lib/j9/system/defs/hostdefs_linux_64.ijs
usr/lib/j9/system/defs/hostdefs_sunos.ijs
usr/lib/j9/system/defs/hostdefs_win.ijs
usr/lib/j9/system/defs/hostdefs_win_64.ijs
usr/lib/j9/system/defs/netdefs_aix.ijs
usr/lib/j9/system/defs/netdefs_android.ijs
usr/lib/j9/system/defs/netdefs_android_64.ijs
usr/lib/j9/system/defs/netdefs_darwin.ijs
usr/lib/j9/system/defs/netdefs_darwin_64.ijs
usr/lib/j9/system/defs/netdefs_linux.ijs
usr/lib/j9/system/defs/netdefs_linux_64.ijs
usr/lib/j9/system/defs/netdefs_sunos.ijs
usr/lib/j9/system/defs/netdefs_win.ijs
usr/lib/j9/system/defs/netdefs_win_64.ijs
usr/lib/j9/system/main/
usr/lib/j9/system/main/regex.ijs
usr/lib/j9/system/main/socket.ijs
usr/lib/j9/system/main/stdlib.ijs
usr/lib/j9/system/main/task.ijs
usr/lib/j9/system/util/
usr/lib/j9/system/util/boot.ijs
usr/lib/j9/system/util/configure.ijs
usr/lib/j9/system/util/edit.png
usr/lib/j9/system/util/edit.svg
usr/lib/j9/system/util/pacman.ijs
usr/lib/j9/system/util/pm.ijs
usr/lib/j9/system/util/pp.ijs
usr/lib/j9/system/util/project.ijs
usr/lib/j9/system/util/scripts.ijs
usr/lib/j9/system/util/tar.ijs
usr/lib/j9/system/util/term.png
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/j9hsd.service
usr/share/
usr/share/applications/
usr/share/applications/j9con.desktop
usr/share/applications/j9hs.desktop
usr/share/applications/j9qt.desktop
usr/share/doc/
usr/share/doc/j9/
usr/share/doc/j9/README.txt
usr/share/doc/j9/jrefcard.pdf
