# files/pixelserv-tls-2.4-1/files
%FILES%
usr/
usr/bin/
usr/bin/pixelserv-ca-init
usr/bin/pixelserv-tls
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pixelserv-tls.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/pixelserv-tls.1.gz
