# files/nordvpn-bin-3.12.0_1-1/files
%FILES%
etc/
usr/
usr/bin/
usr/bin/nordvpn
usr/bin/nordvpnd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nordvpnd.service
usr/lib/systemd/system/nordvpnd.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/nordvpn-bin.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/nordvpn.conf
usr/share/
usr/share/applications/
usr/share/applications/nordvpn.desktop
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/nordvpn
usr/share/doc/
usr/share/doc/nordvpn/
usr/share/doc/nordvpn/changelog.gz
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/nordvpn.1.gz
usr/share/zsh/
usr/share/zsh/functions/
usr/share/zsh/functions/Completion/
usr/share/zsh/functions/Completion/Unix/
usr/share/zsh/functions/Completion/Unix/_nordvpn_auto_complete
var/
var/lib/
var/lib/nordvpn/
var/lib/nordvpn/data/
var/lib/nordvpn/data/cybersec.dat
var/lib/nordvpn/data/insights.dat
var/lib/nordvpn/data/ovpn_template.xslt
var/lib/nordvpn/data/ovpn_xor_template.xslt
var/lib/nordvpn/data/rsa-key-1.pub
var/lib/nordvpn/data/servers.dat
var/lib/nordvpn/icon.svg
var/lib/nordvpn/openvpn
