# files/opensnitch-1.4.3-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/opensnitch
etc/opensnitchd/
etc/opensnitchd/default-config.json
etc/opensnitchd/rules/
etc/opensnitchd/system-fw.json
usr/
usr/bin/
usr/bin/opensnitch-ui
usr/bin/opensnitchd
usr/lib/
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/opensnitch/
usr/lib/python3.9/site-packages/opensnitch/__init__.py
usr/lib/python3.9/site-packages/opensnitch/__pycache__/
usr/lib/python3.9/site-packages/opensnitch/__pycache__/__init__.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/config.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/config.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/customwidgets.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/customwidgets.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/database.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/database.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/desktop_parser.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/desktop_parser.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/nodes.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/nodes.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/resources_rc.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/resources_rc.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/service.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/service.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/ui_pb2.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/ui_pb2.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/ui_pb2_grpc.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/ui_pb2_grpc.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/utils.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/utils.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/version.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/__pycache__/version.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/config.py
usr/lib/python3.9/site-packages/opensnitch/customwidgets.py
usr/lib/python3.9/site-packages/opensnitch/database.py
usr/lib/python3.9/site-packages/opensnitch/desktop_parser.py
usr/lib/python3.9/site-packages/opensnitch/dialogs/
usr/lib/python3.9/site-packages/opensnitch/dialogs/__init__.py
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/__init__.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/preferences.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/preferences.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/processdetails.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/processdetails.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/prompt.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/prompt.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/ruleseditor.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/ruleseditor.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/stats.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/__pycache__/stats.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/dialogs/preferences.py
usr/lib/python3.9/site-packages/opensnitch/dialogs/processdetails.py
usr/lib/python3.9/site-packages/opensnitch/dialogs/prompt.py
usr/lib/python3.9/site-packages/opensnitch/dialogs/ruleseditor.py
usr/lib/python3.9/site-packages/opensnitch/dialogs/stats.py
usr/lib/python3.9/site-packages/opensnitch/nodes.py
usr/lib/python3.9/site-packages/opensnitch/res/
usr/lib/python3.9/site-packages/opensnitch/res/__init__.py
usr/lib/python3.9/site-packages/opensnitch/res/__pycache__/
usr/lib/python3.9/site-packages/opensnitch/res/__pycache__/__init__.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/opensnitch/res/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/opensnitch/res/icon-alert.png
usr/lib/python3.9/site-packages/opensnitch/res/icon-off.png
usr/lib/python3.9/site-packages/opensnitch/res/icon-pause.png
usr/lib/python3.9/site-packages/opensnitch/res/icon-pause.svg
usr/lib/python3.9/site-packages/opensnitch/res/icon-red.png
usr/lib/python3.9/site-packages/opensnitch/res/icon-white.png
usr/lib/python3.9/site-packages/opensnitch/res/icon-white.svg
usr/lib/python3.9/site-packages/opensnitch/res/icon.png
usr/lib/python3.9/site-packages/opensnitch/res/preferences.ui
usr/lib/python3.9/site-packages/opensnitch/res/process_details.ui
usr/lib/python3.9/site-packages/opensnitch/res/prompt.ui
usr/lib/python3.9/site-packages/opensnitch/res/resources.qrc
usr/lib/python3.9/site-packages/opensnitch/res/ruleseditor.ui
usr/lib/python3.9/site-packages/opensnitch/res/stats.ui
usr/lib/python3.9/site-packages/opensnitch/resources_rc.py
usr/lib/python3.9/site-packages/opensnitch/service.py
usr/lib/python3.9/site-packages/opensnitch/ui_pb2.py
usr/lib/python3.9/site-packages/opensnitch/ui_pb2_grpc.py
usr/lib/python3.9/site-packages/opensnitch/utils.py
usr/lib/python3.9/site-packages/opensnitch/version.py
usr/lib/python3.9/site-packages/opensnitch_ui-1.4.2-py3.9.egg-info/
usr/lib/python3.9/site-packages/opensnitch_ui-1.4.2-py3.9.egg-info/PKG-INFO
usr/lib/python3.9/site-packages/opensnitch_ui-1.4.2-py3.9.egg-info/SOURCES.txt
usr/lib/python3.9/site-packages/opensnitch_ui-1.4.2-py3.9.egg-info/dependency_links.txt
usr/lib/python3.9/site-packages/opensnitch_ui-1.4.2-py3.9.egg-info/not-zip-safe
usr/lib/python3.9/site-packages/opensnitch_ui-1.4.2-py3.9.egg-info/top_level.txt
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/opensnitchd.service
usr/share/
usr/share/applications/
usr/share/applications/opensnitch_ui.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/opensnitch-ui.png
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/opensnitch-ui.svg
usr/share/kservices5/
usr/share/kservices5/kcm_opensnitch.desktop
