# files/oddjob-0.34.7-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/oddjob-mkhomedir.conf
etc/dbus-1/system.d/oddjob.conf
etc/oddjob/
etc/oddjobd.conf
etc/oddjobd.conf.d/
etc/oddjobd.conf.d/oddjobd-introspection.conf
etc/oddjobd.conf.d/oddjobd-mkhomedir.conf
usr/
usr/bin/
usr/bin/oddjob_request
usr/bin/oddjobd
usr/lib/
usr/lib/oddjob/
usr/lib/oddjob/oddjob/
usr/lib/oddjob/oddjob/mkhomedir
usr/lib/oddjob/oddjob/sanity.sh
usr/lib/security/
usr/lib/security/pam_oddjob_mkhomedir.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/oddjobd.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/oddjob_request.1.gz
usr/share/man/man5/
usr/share/man/man5/oddjob-mkhomedir.conf.5.gz
usr/share/man/man5/oddjob.conf.5.gz
usr/share/man/man5/oddjobd-introspection.conf.5.gz
usr/share/man/man5/oddjobd-mkhomedir.conf.5.gz
usr/share/man/man5/oddjobd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/oddjobd.8.gz
usr/share/man/man8/pam_oddjob_mkhomedir.8.gz
