# files/metronome-git-3.14.6-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/metronome
etc/metronome/
etc/metronome/certs/
etc/metronome/certs/Makefile
etc/metronome/certs/localhost.cert
etc/metronome/certs/localhost.cnf
etc/metronome/certs/localhost.key
etc/metronome/certs/openssl.cnf
etc/metronome/metronome.cfg.lua
etc/metronome/templates/
etc/metronome/templates/associate.template.txt
etc/metronome/templates/register.template.txt
etc/metronome/templates/reset.template.txt
run/
run/metronome/
usr/
usr/bin/
usr/bin/metronome
usr/bin/metronomectl
usr/lib/
usr/lib/metronome/
usr/lib/metronome/core/
usr/lib/metronome/core/configmanager.lua
usr/lib/metronome/core/hostmanager.lua
usr/lib/metronome/core/loggingmanager.lua
usr/lib/metronome/core/moduleapi.lua
usr/lib/metronome/core/modulemanager.lua
usr/lib/metronome/core/portmanager.lua
usr/lib/metronome/core/sessionmanager.lua
usr/lib/metronome/core/storagemanager.lua
usr/lib/metronome/core/usermanager.lua
usr/lib/metronome/metronome.version
usr/lib/metronome/modules/
usr/lib/metronome/modules/acdf_aux.lib.lua
usr/lib/metronome/modules/adhoc/
usr/lib/metronome/modules/adhoc/adhoc.lib.lua
usr/lib/metronome/modules/adhoc/mod_adhoc.lua
usr/lib/metronome/modules/admin_web/
usr/lib/metronome/modules/admin_web/mod_admin_web.lua
usr/lib/metronome/modules/admin_web/www_files/
usr/lib/metronome/modules/admin_web/www_files/css/
usr/lib/metronome/modules/admin_web/www_files/css/bootstrap-1.4.0.min.css
usr/lib/metronome/modules/admin_web/www_files/css/style.css
usr/lib/metronome/modules/admin_web/www_files/images/
usr/lib/metronome/modules/admin_web/www_files/images/bidi.png
usr/lib/metronome/modules/admin_web/www_files/images/compressed.png
usr/lib/metronome/modules/admin_web/www_files/images/csi-active.png
usr/lib/metronome/modules/admin_web/www_files/images/csi-inactive.png
usr/lib/metronome/modules/admin_web/www_files/images/direct-tls.png
usr/lib/metronome/modules/admin_web/www_files/images/encrypted.png
usr/lib/metronome/modules/admin_web/www_files/images/haze_orange.png
usr/lib/metronome/modules/admin_web/www_files/images/metronome.png
usr/lib/metronome/modules/admin_web/www_files/images/secure.png
usr/lib/metronome/modules/admin_web/www_files/images/sm.png
usr/lib/metronome/modules/admin_web/www_files/index.html
usr/lib/metronome/modules/admin_web/www_files/js/
usr/lib/metronome/modules/admin_web/www_files/js/adhoc.js
usr/lib/metronome/modules/admin_web/www_files/js/jquery-3.5.0.min.js
usr/lib/metronome/modules/admin_web/www_files/js/main.js
usr/lib/metronome/modules/admin_web/www_files/js/strophe.min.js
usr/lib/metronome/modules/favicon/
usr/lib/metronome/modules/favicon/favicon.ico
usr/lib/metronome/modules/favicon/favicon.png
usr/lib/metronome/modules/favicon/mod_favicon.lua
usr/lib/metronome/modules/http_upload_external/
usr/lib/metronome/modules/http_upload_external/mod_http_upload_external.lua
usr/lib/metronome/modules/http_upload_external/resources/
usr/lib/metronome/modules/http_upload_external/resources/share.php
usr/lib/metronome/modules/incidents_handling/
usr/lib/metronome/modules/incidents_handling/incidents_handling.lib.lua
usr/lib/metronome/modules/incidents_handling/mod_incidents_handling.lua
usr/lib/metronome/modules/mam/
usr/lib/metronome/modules/mam/mam.lib.lua
usr/lib/metronome/modules/mam/mod_mam.lua
usr/lib/metronome/modules/mam/validate.lib.lua
usr/lib/metronome/modules/mam_browser/
usr/lib/metronome/modules/mam_browser/mod_mam_browser.lua
usr/lib/metronome/modules/mam_browser/template/
usr/lib/metronome/modules/mam_browser/template/browser.html
usr/lib/metronome/modules/mam_browser/template/css/
usr/lib/metronome/modules/mam_browser/template/css/style.css
usr/lib/metronome/modules/mam_browser/template/fail.html
usr/lib/metronome/modules/mam_browser/template/images/
usr/lib/metronome/modules/mam_browser/template/images/header.png
usr/lib/metronome/modules/mam_browser/template/images/tile.png
usr/lib/metronome/modules/mam_browser/template/login.html
usr/lib/metronome/modules/mam_browser/template/unsecure.html
usr/lib/metronome/modules/mod_acdf.lua
usr/lib/metronome/modules/mod_adhoc_cm.lua
usr/lib/metronome/modules/mod_admin_adhoc.lua
usr/lib/metronome/modules/mod_admin_telnet.lua
usr/lib/metronome/modules/mod_announce.lua
usr/lib/metronome/modules/mod_auth_anonymous.lua
usr/lib/metronome/modules/mod_auth_external.lua
usr/lib/metronome/modules/mod_auth_internal_hashed.lua
usr/lib/metronome/modules/mod_auth_internal_plain.lua
usr/lib/metronome/modules/mod_bidi.lua
usr/lib/metronome/modules/mod_bind_session.lua
usr/lib/metronome/modules/mod_bookmarks.lua
usr/lib/metronome/modules/mod_bosh.lua
usr/lib/metronome/modules/mod_c2s.lua
usr/lib/metronome/modules/mod_component.lua
usr/lib/metronome/modules/mod_compression.lua
usr/lib/metronome/modules/mod_dialback.lua
usr/lib/metronome/modules/mod_disco.lua
usr/lib/metronome/modules/mod_dump_json.lua
usr/lib/metronome/modules/mod_extdisco.lua
usr/lib/metronome/modules/mod_gate_guard.lua
usr/lib/metronome/modules/mod_gdpr.lua
usr/lib/metronome/modules/mod_groups.lua
usr/lib/metronome/modules/mod_http.lua
usr/lib/metronome/modules/mod_http_errors.lua
usr/lib/metronome/modules/mod_http_upload.lua
usr/lib/metronome/modules/mod_iq.lua
usr/lib/metronome/modules/mod_jid_prep.lua
usr/lib/metronome/modules/mod_jingle_nodes.lua
usr/lib/metronome/modules/mod_lastactivity.lua
usr/lib/metronome/modules/mod_message.lua
usr/lib/metronome/modules/mod_message_carbons.lua
usr/lib/metronome/modules/mod_messagefilter.lua
usr/lib/metronome/modules/mod_motd.lua
usr/lib/metronome/modules/mod_muc_bob.lua
usr/lib/metronome/modules/mod_muc_limits.lua
usr/lib/metronome/modules/mod_muc_log.lua
usr/lib/metronome/modules/mod_muc_log_mam.lua
usr/lib/metronome/modules/mod_muc_moderation.lua
usr/lib/metronome/modules/mod_muc_vcard.lua
usr/lib/metronome/modules/mod_net_multiplex.lua
usr/lib/metronome/modules/mod_offline.lua
usr/lib/metronome/modules/mod_pastebin.lua
usr/lib/metronome/modules/mod_pep.lua
usr/lib/metronome/modules/mod_ping.lua
usr/lib/metronome/modules/mod_posix.lua
usr/lib/metronome/modules/mod_presence.lua
usr/lib/metronome/modules/mod_private.lua
usr/lib/metronome/modules/mod_proxy65.lua
usr/lib/metronome/modules/mod_public_service.lua
usr/lib/metronome/modules/mod_pubsub.lua
usr/lib/metronome/modules/mod_push.lua
usr/lib/metronome/modules/mod_register.lua
usr/lib/metronome/modules/mod_register_redirect.lua
usr/lib/metronome/modules/mod_room_zapanddestroy.lua
usr/lib/metronome/modules/mod_roster.lua
usr/lib/metronome/modules/mod_router.lua
usr/lib/metronome/modules/mod_s2s/
usr/lib/metronome/modules/mod_s2s/mod_s2s.lua
usr/lib/metronome/modules/mod_s2s/s2sout.lib.lua
usr/lib/metronome/modules/mod_sasl_s2s.lua
usr/lib/metronome/modules/mod_saslauth.lua
usr/lib/metronome/modules/mod_sec_labels.lua
usr/lib/metronome/modules/mod_server_presence.lua
usr/lib/metronome/modules/mod_server_status.lua
usr/lib/metronome/modules/mod_service_directory.lua
usr/lib/metronome/modules/mod_sic.lua
usr/lib/metronome/modules/mod_stanza_counter.lua
usr/lib/metronome/modules/mod_stanza_optimizations.lua
usr/lib/metronome/modules/mod_storage_cache.lua
usr/lib/metronome/modules/mod_storage_internal.lua
usr/lib/metronome/modules/mod_storage_sql.lua
usr/lib/metronome/modules/mod_stream_management.lua
usr/lib/metronome/modules/mod_subscription_block.lua
usr/lib/metronome/modules/mod_time.lua
usr/lib/metronome/modules/mod_tls.lua
usr/lib/metronome/modules/mod_turnadmin.lua
usr/lib/metronome/modules/mod_uptime.lua
usr/lib/metronome/modules/mod_vcard.lua
usr/lib/metronome/modules/mod_version.lua
usr/lib/metronome/modules/mod_vjud.lua
usr/lib/metronome/modules/mod_watchregistrations.lua
usr/lib/metronome/modules/mod_websocket.lua
usr/lib/metronome/modules/muc/
usr/lib/metronome/modules/muc/mod_muc.lua
usr/lib/metronome/modules/muc/muc.lib.lua
usr/lib/metronome/modules/muc_log_http/
usr/lib/metronome/modules/muc_log_http/generate_log
usr/lib/metronome/modules/muc_log_http/mod_muc_log_http.lua
usr/lib/metronome/modules/muc_log_http/themes/
usr/lib/metronome/modules/muc_log_http/themes/metronome/
usr/lib/metronome/modules/muc_log_http/themes/metronome/components_bit.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/day_body.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/day_dayLink.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/day_message.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/day_messageMe.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/day_time.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/day_titleChange.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/days_bit.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/days_body.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/days_rooms_bit.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/doc.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/month_day.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/month_emptyDay.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/month_footer.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/month_header.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/month_weekDay.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/rooms_bit.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/rooms_body.html
usr/lib/metronome/modules/muc_log_http/themes/metronome/year_title.html
usr/lib/metronome/modules/pep_aux.lib.lua
usr/lib/metronome/modules/privacy/
usr/lib/metronome/modules/privacy/mod_privacy.lua
usr/lib/metronome/modules/privacy/privacy.lib.lua
usr/lib/metronome/modules/pubsub_aux.lib.lua
usr/lib/metronome/modules/register_api/
usr/lib/metronome/modules/register_api/mod_register_api.lua
usr/lib/metronome/modules/register_api/send_mail
usr/lib/metronome/modules/register_api/template/
usr/lib/metronome/modules/register_api/template/associate_fail_t.html
usr/lib/metronome/modules/register_api/template/associate_form_t.html
usr/lib/metronome/modules/register_api/template/associate_success_t.html
usr/lib/metronome/modules/register_api/template/css/
usr/lib/metronome/modules/register_api/template/css/style.css
usr/lib/metronome/modules/register_api/template/images/
usr/lib/metronome/modules/register_api/template/images/header.png
usr/lib/metronome/modules/register_api/template/images/tile.png
usr/lib/metronome/modules/register_api/template/reset_fail_t.html
usr/lib/metronome/modules/register_api/template/reset_form_t.html
usr/lib/metronome/modules/register_api/template/reset_nomatch_t.html
usr/lib/metronome/modules/register_api/template/reset_password_check_t.html
usr/lib/metronome/modules/register_api/template/reset_success_t.html
usr/lib/metronome/modules/register_api/template/verify_fail_t.html
usr/lib/metronome/modules/register_api/template/verify_form_t.html
usr/lib/metronome/modules/register_api/template/verify_success_t.html
usr/lib/metronome/modules/sasl_aux.lib.lua
usr/lib/metronome/modules/spim_block/
usr/lib/metronome/modules/spim_block/mod_spim_block.lua
usr/lib/metronome/modules/spim_block/template/
usr/lib/metronome/modules/spim_block/template/css/
usr/lib/metronome/modules/spim_block/template/css/style.css
usr/lib/metronome/modules/spim_block/template/fail.html
usr/lib/metronome/modules/spim_block/template/form.html
usr/lib/metronome/modules/spim_block/template/images/
usr/lib/metronome/modules/spim_block/template/images/header.png
usr/lib/metronome/modules/spim_block/template/images/tile.png
usr/lib/metronome/modules/spim_block/template/verify.html
usr/lib/metronome/net/
usr/lib/metronome/net/adns.lua
usr/lib/metronome/net/dns.lua
usr/lib/metronome/net/http.lua
usr/lib/metronome/net/http/
usr/lib/metronome/net/http/codes.lua
usr/lib/metronome/net/http/parser.lua
usr/lib/metronome/net/http/server.lua
usr/lib/metronome/net/server.lua
usr/lib/metronome/net/server_event.lua
usr/lib/metronome/util/
usr/lib/metronome/util/address_selection.lua
usr/lib/metronome/util/array.lua
usr/lib/metronome/util/auxiliary.lua
usr/lib/metronome/util/caps.lua
usr/lib/metronome/util/certmanager.lua
usr/lib/metronome/util/dataforms.lua
usr/lib/metronome/util/datamanager.lua
usr/lib/metronome/util/datetime.lua
usr/lib/metronome/util/debug.lua
usr/lib/metronome/util/dependencies.lua
usr/lib/metronome/util/encodings.so
usr/lib/metronome/util/envload.lua
usr/lib/metronome/util/events.lua
usr/lib/metronome/util/filters.lua
usr/lib/metronome/util/hashes.so
usr/lib/metronome/util/helpers.lua
usr/lib/metronome/util/hmac.lua
usr/lib/metronome/util/httpstream.lua
usr/lib/metronome/util/import.lua
usr/lib/metronome/util/ip.lua
usr/lib/metronome/util/iterators.lua
usr/lib/metronome/util/jid.lua
usr/lib/metronome/util/json.lua
usr/lib/metronome/util/jsonload.lua
usr/lib/metronome/util/logger.lua
usr/lib/metronome/util/mail.lua
usr/lib/metronome/util/metronomectl.lua
usr/lib/metronome/util/multitable.lua
usr/lib/metronome/util/openssl.lua
usr/lib/metronome/util/pluginloader.lua
usr/lib/metronome/util/pposix.so
usr/lib/metronome/util/pubsub.lua
usr/lib/metronome/util/rostermanager.lua
usr/lib/metronome/util/s2smanager.lua
usr/lib/metronome/util/sasl.lua
usr/lib/metronome/util/sasl/
usr/lib/metronome/util/sasl/anonymous.lua
usr/lib/metronome/util/sasl/digest-md5.lua
usr/lib/metronome/util/sasl/external.lua
usr/lib/metronome/util/sasl/plain.lua
usr/lib/metronome/util/sasl/scram.lua
usr/lib/metronome/util/serialization.lua
usr/lib/metronome/util/set.lua
usr/lib/metronome/util/signal.so
usr/lib/metronome/util/stanza.lua
usr/lib/metronome/util/termcolours.lua
usr/lib/metronome/util/throttle.lua
usr/lib/metronome/util/timer.lua
usr/lib/metronome/util/uuid.lua
usr/lib/metronome/util/watchdog.lua
usr/lib/metronome/util/websocket.lua
usr/lib/metronome/util/x509.lua
usr/lib/metronome/util/xmppstream.lua
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/metronome.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/metronome-git.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/metronome-git.conf
usr/share/
usr/share/licences/
usr/share/licences/metronome-git/
usr/share/licences/metronome-git/LICENCE
var/
var/lib/
var/lib/metronome/
var/log/
var/log/metronome/
