# files/openvpn3-13_beta-2/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/net.openvpn.v3.client.conf
etc/dbus-1/system.d/net.openvpn.v3.conf
etc/dbus-1/system.d/net.openvpn.v3.configuration.conf
etc/dbus-1/system.d/net.openvpn.v3.log.conf
etc/dbus-1/system.d/net.openvpn.v3.netcfg.conf
etc/dbus-1/system.d/net.openvpn.v3.sessions.conf
etc/openvpn3/
etc/openvpn3/autoload/
usr/
usr/bin/
usr/bin/openvpn2
usr/bin/openvpn3
usr/bin/openvpn3-admin
usr/bin/openvpn3-as
usr/bin/openvpn3-autoload
usr/lib/
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/openvpn3/
usr/lib/python3.9/site-packages/openvpn3/ConfigManager.py
usr/lib/python3.9/site-packages/openvpn3/ConfigParser.py
usr/lib/python3.9/site-packages/openvpn3/NetCfgManager.py
usr/lib/python3.9/site-packages/openvpn3/SessionManager.py
usr/lib/python3.9/site-packages/openvpn3/__init__.py
usr/lib/python3.9/site-packages/openvpn3/__pycache__/
usr/lib/python3.9/site-packages/openvpn3/__pycache__/ConfigManager.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/ConfigManager.cpython-39.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/ConfigParser.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/ConfigParser.cpython-39.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/NetCfgManager.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/NetCfgManager.cpython-39.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/SessionManager.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/SessionManager.cpython-39.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/__init__.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/constants.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/openvpn3/__pycache__/constants.cpython-39.pyc
usr/lib/python3.9/site-packages/openvpn3/constants.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/openvpn3-autoload.service
usr/libexec/
usr/libexec/openvpn3-linux/
usr/libexec/openvpn3-linux/openvpn3-service-backendstart
usr/libexec/openvpn3-linux/openvpn3-service-client
usr/libexec/openvpn3-linux/openvpn3-service-configmgr
usr/libexec/openvpn3-linux/openvpn3-service-logger
usr/libexec/openvpn3-linux/openvpn3-service-netcfg
usr/libexec/openvpn3-linux/openvpn3-service-sessionmgr
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/openvpn2
usr/share/bash-completion/completions/openvpn3
usr/share/bash-completion/completions/openvpn3-admin
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/net.openvpn.v3.backends.service
usr/share/dbus-1/system-services/net.openvpn.v3.configuration.service
usr/share/dbus-1/system-services/net.openvpn.v3.log.service
usr/share/dbus-1/system-services/net.openvpn.v3.netcfg.service
usr/share/dbus-1/system-services/net.openvpn.v3.sessions.service
usr/share/polkit-1/
usr/share/polkit-1/rules.d/
usr/share/polkit-1/rules.d/net.openvpn.v3.rules
var/
var/lib/
var/lib/openvpn3/
var/lib/openvpn3/configs/
