# files/uresourced-0.5.3-1/files
%FILES%
etc/
etc/uresourced.conf
usr/
usr/lib/
usr/lib/cgroupify
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/uresourced.service
usr/lib/systemd/system/user@.service.d/
usr/lib/systemd/system/user@.service.d/00-uresourced.conf
usr/lib/systemd/user/
usr/lib/systemd/user/app-chromium\\x2dbrowser-.scope.d/
usr/lib/systemd/user/app-chromium\\x2dbrowser-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-chromium\\x2dbrowser@@.service.d/
usr/lib/systemd/user/app-chromium\\x2dbrowser@@.service.d/99-cgroupify.conf
usr/lib/systemd/user/app-firefox-.scope.d/
usr/lib/systemd/user/app-firefox-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-firefox@@.service.d/
usr/lib/systemd/user/app-firefox@@.service.d/99-cgroupify.conf
usr/lib/systemd/user/app-glib-chromium\\x2dbrowser-.scope.d/
usr/lib/systemd/user/app-glib-chromium\\x2dbrowser-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-glib-firefox-.scope.d/
usr/lib/systemd/user/app-glib-firefox-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-glib-google\\x2dchrome-.scope.d/
usr/lib/systemd/user/app-glib-google\\x2dchrome-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-glib-org.gnome.Epiphany-.scope.d/
usr/lib/systemd/user/app-glib-org.gnome.Epiphany-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-gnome-chromium\\x2dbrowser-.scope.d/
usr/lib/systemd/user/app-gnome-chromium\\x2dbrowser-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-gnome-firefox-.scope.d/
usr/lib/systemd/user/app-gnome-firefox-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-gnome-google\\x2dchrome-.scope.d/
usr/lib/systemd/user/app-gnome-google\\x2dchrome-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-gnome-org.gnome.Epiphany-.scope.d/
usr/lib/systemd/user/app-gnome-org.gnome.Epiphany-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-google\\x2dchrome-.scope.d/
usr/lib/systemd/user/app-google\\x2dchrome-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-google\\x2dchrome@@.service.d/
usr/lib/systemd/user/app-google\\x2dchrome@@.service.d/99-cgroupify.conf
usr/lib/systemd/user/app-org.gnome.Epiphany-.scope.d/
usr/lib/systemd/user/app-org.gnome.Epiphany-.scope.d/99-cgroupify.conf
usr/lib/systemd/user/app-org.gnome.Epiphany@@.service.d/
usr/lib/systemd/user/app-org.gnome.Epiphany@@.service.d/99-cgroupify.conf
usr/lib/systemd/user/at-spi-dbus-bus.service.d/
usr/lib/systemd/user/at-spi-dbus-bus.service.d/00-uresourced.conf
usr/lib/systemd/user/cgroupify@.service
usr/lib/systemd/user/dbus-.service.d/
usr/lib/systemd/user/dbus-.service.d/00-uresourced.conf
usr/lib/systemd/user/dbus-broker.service.d/
usr/lib/systemd/user/dbus-broker.service.d/00-uresourced.conf
usr/lib/systemd/user/evolution-.service.d/
usr/lib/systemd/user/evolution-.service.d/00-uresourced.conf
usr/lib/systemd/user/flatpak-.scope.d/
usr/lib/systemd/user/flatpak-.scope.d/00-uresourced.conf
usr/lib/systemd/user/gnome-launched-.scope.d/
usr/lib/systemd/user/gnome-launched-.scope.d/00-uresourced.conf
usr/lib/systemd/user/gnome-launched-gnome-system-monitor.desktop-.scope.d/
usr/lib/systemd/user/gnome-launched-gnome-system-monitor.desktop-.scope.d/00-uresourced.conf
usr/lib/systemd/user/gnome-launched-org.gnome.Usage.desktop-.scope.d/
usr/lib/systemd/user/gnome-launched-org.gnome.Usage.desktop-.scope.d/00-uresourced.conf
usr/lib/systemd/user/gnome-shell-.service.d/
usr/lib/systemd/user/gnome-shell-.service.d/00-uresourced.conf
usr/lib/systemd/user/graphical-session.target.wants/
usr/lib/systemd/user/graphical-session.target.wants/uresourced.service
usr/lib/systemd/user/gsd-.service.d/
usr/lib/systemd/user/gsd-.service.d/00-uresourced.conf
usr/lib/systemd/user/gvfs-.service.d/
usr/lib/systemd/user/gvfs-.service.d/00-uresourced.conf
usr/lib/systemd/user/pipewire.service.d/
usr/lib/systemd/user/pipewire.service.d/00-uresourced.conf
usr/lib/systemd/user/pulseaudio.service.d/
usr/lib/systemd/user/pulseaudio.service.d/00-uresourced.conf
usr/lib/systemd/user/tracker-.service.d/
usr/lib/systemd/user/tracker-.service.d/00-uresourced.conf
usr/lib/systemd/user/uresourced.service
usr/lib/systemd/user/zeitgeist-.service.d/
usr/lib/systemd/user/zeitgeist-.service.d/00-uresourced.conf
usr/lib/systemd/user/zeitgeist.service.d/
usr/lib/systemd/user/zeitgeist.service.d/00-uresourced.conf
usr/lib/uresourced
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.freedesktop.UResourced.conf
usr/share/doc/
usr/share/doc/uresourced/
usr/share/doc/uresourced/NEWS.md
usr/share/doc/uresourced/README
