# files/systemd-guest-user-33.2-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/home-guest.mount
usr/lib/systemd/system/sysinit.target.wants/
usr/lib/systemd/system/sysinit.target.wants/systemd-guest-config.service
usr/lib/systemd/system/sysinit.target.wants/systemd-guest-user.service
usr/lib/systemd/system/systemd-guest-config.service
usr/lib/systemd/system/systemd-guest-user.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/guest.conf
usr/share/
usr/share/licenses/
usr/share/licenses/systemd-guest-user/
usr/share/licenses/systemd-guest-user/LICENSE.md
