# files/naiveproxy-103.0.5060.53_3-1/files
%FILES%
etc/
etc/naiveproxy/
etc/naiveproxy/config.json
usr/
usr/bin/
usr/bin/naiveproxy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/naiveproxy.service
usr/lib/systemd/system/naiveproxy@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/naiveproxy.conf
usr/share/
usr/share/doc/
usr/share/doc/naiveproxy/
usr/share/doc/naiveproxy/README.md
usr/share/doc/naiveproxy/USAGE.txt
usr/share/licenses/
usr/share/licenses/naiveproxy/
usr/share/licenses/naiveproxy/LICENSE
