# files/cpupower-gui-git-1.0.0.r14.g6c6e8fa-1/files
%FILES%
etc/
etc/cpupower_gui.conf
etc/cpupower_gui.d/
etc/cpupower_gui.d/README
etc/cpupower_gui.d/my_profile.profile.ex
usr/
usr/bin/
usr/bin/cpupower-gui
usr/lib/
usr/lib/cpupower-gui/
usr/lib/cpupower-gui/cpupower-gui-helper
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/cpupower-gui-helper.service
usr/lib/systemd/system/cpupower-gui.service
usr/lib/systemd/user/
usr/lib/systemd/user/cpupower-gui-user.service
usr/share/
usr/share/applications/
usr/share/applications/org.rnd2.cpupower_gui.desktop
usr/share/cpupower-gui/
usr/share/cpupower-gui/cpupower-gui.gresource
usr/share/cpupower-gui/cpupower_gui/
usr/share/cpupower-gui/cpupower_gui/__init__.py
usr/share/cpupower-gui/cpupower_gui/config.py
usr/share/cpupower-gui/cpupower_gui/helper.py
usr/share/cpupower-gui/cpupower_gui/main.py
usr/share/cpupower-gui/cpupower_gui/utils.py
usr/share/cpupower-gui/cpupower_gui/window.py
usr/share/cpupower-gui/scaling-governors.conf
usr/share/dbus-1/
usr/share/dbus-1/services/
usr/share/dbus-1/services/org.rnd2.cpupower_gui.service
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.rnd2.cpupower_gui.helper.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.rnd2.cpupower_gui.helper.conf
usr/share/glib-2.0/
usr/share/glib-2.0/schemas/
usr/share/glib-2.0/schemas/org.rnd2.cpupower_gui.gschema.xml
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/org.rnd2.cpupower_gui.svg
usr/share/locale/
usr/share/locale/el_GR/
usr/share/locale/el_GR/LC_MESSAGES/
usr/share/locale/el_GR/LC_MESSAGES/cpupower-gui.mo
usr/share/locale/en/
usr/share/locale/en/LC_MESSAGES/
usr/share/locale/en/LC_MESSAGES/cpupower-gui.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/cpupower-gui.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/cpupower-gui.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/cpupower-gui.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/cpupower-gui.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/cpupower-gui.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/cpupower-gui.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/cpupower-gui.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/cpupower-gui.1.gz
usr/share/man/zh_CN/
usr/share/man/zh_CN/man1/
usr/share/man/zh_CN/man1/cpupower-gui.1.gz
usr/share/metainfo/
usr/share/metainfo/org.rnd2.cpupower_gui.appdata.xml
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/org.rnd2.cpupower-gui.policy
usr/share/polkit-1/rules.d/
usr/share/polkit-1/rules.d/org.rnd2.cpupower-gui.rules
