# files/miredo-1.2.6+r159+ge5f5652-1/files
%FILES%
etc/
etc/miredo/
etc/miredo/client-hook
etc/miredo/miredo.conf
usr/
usr/bin/
usr/bin/miredo
usr/bin/miredo-checkconf
usr/bin/miredo-server
usr/bin/teredo-mire
usr/include/
usr/include/libteredo/
usr/include/libteredo/teredo-udp.h
usr/include/libteredo/teredo.h
usr/include/libteredo/tunnel.h
usr/include/libtun6/
usr/include/libtun6/tun6.h
usr/lib/
usr/lib/libteredo.so
usr/lib/libteredo.so.6
usr/lib/libteredo.so.6.0.0
usr/lib/libtun6.so
usr/lib/libtun6.so.2
usr/lib/libtun6.so.2.0.0
usr/lib/miredo/
usr/lib/miredo/miredo-privproc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/miredo-server.service
usr/lib/systemd/system/miredo.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/miredo.conf
usr/share/
usr/share/doc/
usr/share/doc/miredo/
usr/share/doc/miredo/examples/
usr/share/doc/miredo/examples/miredo-server.conf
usr/share/doc/miredo/examples/miredo.conf
usr/share/locale/
usr/share/locale/en/
usr/share/locale/en/LC_MESSAGES/
usr/share/locale/en/LC_MESSAGES/miredo.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/miredo.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/miredo.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/miredo.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/teredo-mire.1.gz
usr/share/man/man5/
usr/share/man/man5/miredo-server.conf.5.gz
usr/share/man/man5/miredo.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/miredo-checkconf.8.gz
usr/share/man/man8/miredo-server.8.gz
usr/share/man/man8/miredo.8.gz
