# files/onedrive-abraunegg-2.4.14-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/onedrive
usr/
usr/bin/
usr/bin/onedrive
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/onedrive@.service
usr/lib/systemd/user/
usr/lib/systemd/user/onedrive.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/onedrive
usr/share/doc/
usr/share/doc/onedrive/
usr/share/doc/onedrive/BusinessSharedFolders.md
usr/share/doc/onedrive/CHANGELOG.md
usr/share/doc/onedrive/Docker.md
usr/share/doc/onedrive/INSTALL.md
usr/share/doc/onedrive/LICENSE
usr/share/doc/onedrive/README.md
usr/share/doc/onedrive/SharePoint-Shared-Libraries.md
usr/share/doc/onedrive/USAGE.md
usr/share/doc/onedrive/advanced-usage.md
usr/share/doc/onedrive/config
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/onedrive.fish
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/onedrive.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_onedrive
