# files/jdownloader2-latest-17/files
%FILES%
opt/
opt/JDownloader/
opt/JDownloaderScripts/
opt/JDownloaderScripts/JDownloader
opt/JDownloaderScripts/JDownloaderHeadless
opt/JDownloaderScripts/JDownloaderHeadlessCleanLogin
opt/JDownloaderScripts/JDownloaderHeadlessCtl
opt/JDownloaderScripts/functions.sh
usr/
usr/bin/
usr/bin/JDownloader
usr/bin/JDownloaderHeadless
usr/bin/JDownloaderHeadlessCleanLogin
usr/bin/JDownloaderHeadlessCtl
usr/bin/jdownloader
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/jdownloader.service
usr/share/
usr/share/applications/
usr/share/applications/jd-containers.desktop
usr/share/applications/jdownloader.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/16x16/
usr/share/icons/hicolor/16x16/apps/
usr/share/icons/hicolor/16x16/apps/jdownloader.png
usr/share/icons/hicolor/16x16/mimetypes/
usr/share/icons/hicolor/16x16/mimetypes/jd-container.png
usr/share/icons/hicolor/22x22/
usr/share/icons/hicolor/22x22/apps/
usr/share/icons/hicolor/22x22/apps/jdownloader.png
usr/share/icons/hicolor/22x22/mimetypes/
usr/share/icons/hicolor/22x22/mimetypes/jd-container.png
usr/share/icons/hicolor/24x24/
usr/share/icons/hicolor/24x24/apps/
usr/share/icons/hicolor/24x24/apps/jdownloader.png
usr/share/icons/hicolor/24x24/mimetypes/
usr/share/icons/hicolor/24x24/mimetypes/jd-container.png
usr/share/icons/hicolor/256x256/
usr/share/icons/hicolor/256x256/apps/
usr/share/icons/hicolor/256x256/apps/jdownloader.png
usr/share/icons/hicolor/256x256/mimetypes/
usr/share/icons/hicolor/256x256/mimetypes/jd-container.png
usr/share/icons/hicolor/32x32/
usr/share/icons/hicolor/32x32/apps/
usr/share/icons/hicolor/32x32/apps/jdownloader.png
usr/share/icons/hicolor/32x32/mimetypes/
usr/share/icons/hicolor/32x32/mimetypes/jd-container.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/jdownloader.png
usr/share/icons/hicolor/48x48/mimetypes/
usr/share/icons/hicolor/48x48/mimetypes/jd-container.png
usr/share/mime/
usr/share/mime/packages/
usr/share/mime/packages/jdownloader.xml
