# files/bluetooth-autoconnect-1.3-1.1/files
%FILES%
usr/
usr/bin/
usr/bin/bluetooth-autoconnect
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/bluetooth-autoconnect.service
usr/lib/systemd/user/
usr/lib/systemd/user/pulseaudio-bluetooth-autoconnect.service
usr/share/
usr/share/licenses/
usr/share/licenses/bluetooth-autoconnect/
usr/share/licenses/bluetooth-autoconnect/LICENSE
