# files/garuda-common-settings-git-2.3.0.r0.gcd129aa-1/files
%FILES%
etc/
etc/grub.d/
etc/grub.d/61_custom_leave_options
etc/profile.d/
etc/profile.d/canberra-gtk-module.sh
etc/profile.d/qt5-accessibility.sh
etc/skel/
etc/skel/.config/
etc/skel/.config/autostart/
etc/skel/.config/autostart/conky.desktop
etc/skel/.config/autostart/garuda-welcome.desktop
etc/skel/.config/autostart/initial-user-setup.desktop
etc/skel/.config/autostart/org.flameshot.Flameshot.desktop
etc/skel/.config/autostart/org.ksnip.ksnip.desktop
etc/skel/.config/autostart/redshift-gtk.desktop
etc/skel/.config/autostart/variety.desktop
etc/skel/.config/autostart/xdman.desktop
etc/skel/.config/bleachbit/
etc/skel/.config/bleachbit/bleachbit.ini
etc/skel/.config/deluge/
etc/skel/.config/deluge/gtk3ui.conf
etc/skel/.config/ksnip/
etc/skel/.config/ksnip/ksnip.conf
etc/skel/.config/libreoffice/
etc/skel/.config/libreoffice/4/
etc/skel/.config/libreoffice/4/user/
etc/skel/.config/libreoffice/4/user/registrymodifications.xcu
etc/skel/.config/micro/
etc/skel/.config/micro/plug/
etc/skel/.config/micro/plug/aspell/
etc/skel/.config/micro/plug/aspell/LICENSE
etc/skel/.config/micro/plug/aspell/README.md
etc/skel/.config/micro/plug/aspell/aspell.lua
etc/skel/.config/micro/plug/aspell/help/
etc/skel/.config/micro/plug/aspell/help/aspell.md
etc/skel/.config/micro/plug/aspell/repo.json
etc/skel/.config/micro/plug/fzf/
etc/skel/.config/micro/plug/fzf/README.md
etc/skel/.config/micro/plug/fzf/main.lua
etc/skel/.config/micro/plug/fzf/repo.json
etc/skel/.config/micro/plug/misspell/
etc/skel/.config/micro/plug/misspell/README.md
etc/skel/.config/micro/plug/misspell/misspell.lua
etc/skel/.config/micro/plug/misspell/repo.json
etc/skel/.config/micro/plug/wc/
etc/skel/.config/micro/plug/wc/CHANGELOG.md
etc/skel/.config/micro/plug/wc/LICENSE
etc/skel/.config/micro/plug/wc/README.md
etc/skel/.config/micro/plug/wc/help/
etc/skel/.config/micro/plug/wc/help/wc.md
etc/skel/.config/micro/plug/wc/info.json
etc/skel/.config/micro/plug/wc/repo.json
etc/skel/.config/micro/plug/wc/wc.lua
etc/skel/.config/micro/settings.json
etc/skel/.config/mpv/
etc/skel/.config/mpv/mpv.conf
etc/skel/.config/psd/
etc/skel/.config/psd/psd.conf
etc/skel/.config/pulseaudio-ctl/
etc/skel/.config/pulseaudio-ctl/config
etc/skel/.config/smplayer/
etc/skel/.config/smplayer/smplayer.ini
etc/skel/.config/strawberry/
etc/skel/.config/strawberry/strawberry.conf
etc/skel/.config/variety/
etc/skel/.config/variety/variety.conf
etc/skel/.config/vlc/
etc/skel/.config/vlc/vlcrc
etc/sudoers.d/
etc/sudoers.d/insults
etc/sudoers.d/pwfeedback
usr/
usr/lib/
usr/lib/modprobe.d/
usr/lib/modprobe.d/nobeep.conf
usr/lib/modprobe.d/noime.conf
usr/lib/sysctl.d/
usr/lib/sysctl.d/99-sysctl-garuda.conf
usr/lib/systemd/
usr/lib/systemd/system.conf.d/
usr/lib/systemd/system.conf.d/limits.conf
usr/lib/systemd/system/
usr/lib/systemd/system/garuda-pacman-lock.service
usr/lib/systemd/user.conf.d/
usr/lib/systemd/user.conf.d/limits.conf
usr/lib/systemd/zram-generator.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/40-hpet-permissions.rules
usr/lib/udev/rules.d/50-sata.rules
usr/lib/udev/rules.d/60-ioschedulers.rules
usr/share/
usr/share/applications/
usr/share/applications/reboot-to-uefi.desktop
usr/share/applications/suspend-then-hibernate.desktop
usr/share/polkit-1/
usr/share/polkit-1/rules.d/
usr/share/polkit-1/rules.d/99-garuda.rules
