# files/amdfand-bin-1.0.13-4/files
%FILES%
usr/
usr/bin/
usr/bin/amdfand
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/amdfand.service
