# files/rqlite-7.7.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/rqbench
usr/bin/rqlite
usr/bin/rqlited
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rqlite.service
usr/lib/systemd/system/rqlite@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/rqlite.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/rqlite.conf
