# files/beanstalkd-1.12-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/beanstalkd
usr/
usr/bin/
usr/bin/beanstalkd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/beanstalkd.service
usr/lib/systemd/system/beanstalkd.socket
usr/lib/systemd/system/beanstalkd@.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/beanstalkd.1.gz
