# files/memavaild-0.6-1/files
%FILES%
etc/
etc/memavaild.conf
usr/
usr/bin/
usr/bin/memavaild
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/memavaild.service
usr/share/
usr/share/doc/
usr/share/doc/memavaild/
usr/share/doc/memavaild/MANPAGE.md
usr/share/doc/memavaild/README.md
usr/share/licenses/
usr/share/licenses/memavaild/
usr/share/licenses/memavaild/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/memavaild.8.gz
usr/share/memavaild/
usr/share/memavaild/memavaild.conf
