# files/precached-git-1.7.3-0/files
%FILES%
etc/
etc/precached/
etc/precached/log4rs.yaml
etc/precached/precached.conf
etc/precached/rules.d/
etc/precached/rules.d/00-log-fork-bombs.rules
etc/precached/rules.d/10-cache-on-login.rules
etc/precached/rules.d/99-ping-logger.rules
etc/precached/rules.d/README
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/precached-trigger.desktop
usr/
usr/bin/
usr/bin/iotracectl
usr/bin/precached
usr/bin/precached-debug
usr/bin/precached-trigger
usr/bin/precachedctl
usr/bin/precachedtop
usr/bin/rulesctl
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system-preset/
usr/lib/systemd/system-preset/50-precached.preset
usr/lib/systemd/system/
usr/lib/systemd/system/precached-prime-caches.service
usr/lib/systemd/system/precached-prime-caches.timer
usr/lib/systemd/system/precached.service
usr/lib/systemd/user-preset/
usr/lib/systemd/user-preset/50-precached.preset
usr/lib/systemd/user/
usr/lib/systemd/user/precached-trigger.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/iotracectl
usr/share/bash-completion/completions/precached-debug
usr/share/bash-completion/completions/precached-trigger
usr/share/bash-completion/completions/precachedctl
usr/share/bash-completion/completions/precachedtop
usr/share/bash-completion/completions/rulesctl
usr/share/doc/
usr/share/doc/precached/
usr/share/doc/precached/examples/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/precached-trigger.1.gz
usr/share/man/man1/precachedtop.1.gz
usr/share/man/man5/
usr/share/man/man5/precached.conf.5.gz
usr/share/man/man5/precached.rules.5.gz
usr/share/man/man8/
usr/share/man/man8/iotracectl.8.gz
usr/share/man/man8/precached-debug.8.gz
usr/share/man/man8/precached.8.gz
usr/share/man/man8/precachedctl.8.gz
usr/share/man/man8/rulesctl.8.gz
usr/share/precached/
usr/share/precached/i18n/
usr/share/precached/i18n/c
usr/share/precached/i18n/de_at
usr/share/precached/i18n/de_at.utf-8
usr/share/precached/i18n/de_at.utf8
usr/share/precached/i18n/de_de.utf-8
usr/share/precached/i18n/de_de.utf8
usr/share/precached/i18n/de_de/
usr/share/precached/i18n/de_de/messages.fluent
usr/share/precached/i18n/en_uk
usr/share/precached/i18n/en_uk.utf-8
usr/share/precached/i18n/en_uk.utf8
usr/share/precached/i18n/en_us.utf-8
usr/share/precached/i18n/en_us.utf8
usr/share/precached/i18n/en_us/
usr/share/precached/i18n/en_us/messages.fluent
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_iotracectl
usr/share/zsh/site-functions/_precached-debug
usr/share/zsh/site-functions/_precached-trigger
usr/share/zsh/site-functions/_precachedctl
usr/share/zsh/site-functions/_precachedtop
usr/share/zsh/site-functions/_rulesctl
var/
var/lib/
var/lib/precached/
var/lib/precached/iotrace/
