# files/inspircd-3.10.0-1/files
%FILES%
etc/
etc/inspircd/
etc/inspircd/help.txt
usr/
usr/bin/
usr/bin/inspircd
usr/bin/inspircd-genssl
usr/bin/inspircd-testssl
usr/lib/
usr/lib/inspircd/
usr/lib/inspircd/apparmor
usr/lib/inspircd/inspircd
usr/lib/inspircd/inspircd.service
usr/lib/inspircd/logrotate
usr/lib/inspircd/manuals/
usr/lib/inspircd/manuals/inspircd-genssl.1
usr/lib/inspircd/manuals/inspircd-testssl.1
usr/lib/inspircd/manuals/inspircd.1
usr/lib/inspircd/modules/
usr/lib/inspircd/modules/core_channel.so
usr/lib/inspircd/modules/core_dns.so
usr/lib/inspircd/modules/core_hostname_lookup.so
usr/lib/inspircd/modules/core_info.so
usr/lib/inspircd/modules/core_list.so
usr/lib/inspircd/modules/core_loadmodule.so
usr/lib/inspircd/modules/core_lusers.so
usr/lib/inspircd/modules/core_message.so
usr/lib/inspircd/modules/core_mode.so
usr/lib/inspircd/modules/core_oper.so
usr/lib/inspircd/modules/core_reloadmodule.so
usr/lib/inspircd/modules/core_serialize_rfc.so
usr/lib/inspircd/modules/core_stats.so
usr/lib/inspircd/modules/core_stub.so
usr/lib/inspircd/modules/core_user.so
usr/lib/inspircd/modules/core_wallops.so
usr/lib/inspircd/modules/core_who.so
usr/lib/inspircd/modules/core_whois.so
usr/lib/inspircd/modules/core_whowas.so
usr/lib/inspircd/modules/core_xline.so
usr/lib/inspircd/modules/m_abbreviation.so
usr/lib/inspircd/modules/m_alias.so
usr/lib/inspircd/modules/m_allowinvite.so
usr/lib/inspircd/modules/m_alltime.so
usr/lib/inspircd/modules/m_anticaps.so
usr/lib/inspircd/modules/m_argon2.so
usr/lib/inspircd/modules/m_auditorium.so
usr/lib/inspircd/modules/m_autoop.so
usr/lib/inspircd/modules/m_banexception.so
usr/lib/inspircd/modules/m_banredirect.so
usr/lib/inspircd/modules/m_bcrypt.so
usr/lib/inspircd/modules/m_blockamsg.so
usr/lib/inspircd/modules/m_blockcaps.so
usr/lib/inspircd/modules/m_blockcolor.so
usr/lib/inspircd/modules/m_botmode.so
usr/lib/inspircd/modules/m_callerid.so
usr/lib/inspircd/modules/m_cap.so
usr/lib/inspircd/modules/m_cban.so
usr/lib/inspircd/modules/m_censor.so
usr/lib/inspircd/modules/m_cgiirc.so
usr/lib/inspircd/modules/m_chancreate.so
usr/lib/inspircd/modules/m_chanfilter.so
usr/lib/inspircd/modules/m_chanhistory.so
usr/lib/inspircd/modules/m_chanlog.so
usr/lib/inspircd/modules/m_channames.so
usr/lib/inspircd/modules/m_channelban.so
usr/lib/inspircd/modules/m_check.so
usr/lib/inspircd/modules/m_chghost.so
usr/lib/inspircd/modules/m_chgident.so
usr/lib/inspircd/modules/m_chgname.so
usr/lib/inspircd/modules/m_classban.so
usr/lib/inspircd/modules/m_clearchan.so
usr/lib/inspircd/modules/m_cloaking.so
usr/lib/inspircd/modules/m_clones.so
usr/lib/inspircd/modules/m_codepage.so
usr/lib/inspircd/modules/m_commonchans.so
usr/lib/inspircd/modules/m_conn_join.so
usr/lib/inspircd/modules/m_conn_umodes.so
usr/lib/inspircd/modules/m_conn_waitpong.so
usr/lib/inspircd/modules/m_connectban.so
usr/lib/inspircd/modules/m_connflood.so
usr/lib/inspircd/modules/m_customprefix.so
usr/lib/inspircd/modules/m_customtitle.so
usr/lib/inspircd/modules/m_cycle.so
usr/lib/inspircd/modules/m_dccallow.so
usr/lib/inspircd/modules/m_deaf.so
usr/lib/inspircd/modules/m_delayjoin.so
usr/lib/inspircd/modules/m_delaymsg.so
usr/lib/inspircd/modules/m_denychans.so
usr/lib/inspircd/modules/m_disable.so
usr/lib/inspircd/modules/m_dnsbl.so
usr/lib/inspircd/modules/m_exemptchanops.so
usr/lib/inspircd/modules/m_filter.so
usr/lib/inspircd/modules/m_flashpolicyd.so
usr/lib/inspircd/modules/m_gecosban.so
usr/lib/inspircd/modules/m_geo_maxmind.so
usr/lib/inspircd/modules/m_geoban.so
usr/lib/inspircd/modules/m_geoclass.so
usr/lib/inspircd/modules/m_globalload.so
usr/lib/inspircd/modules/m_globops.so
usr/lib/inspircd/modules/m_haproxy.so
usr/lib/inspircd/modules/m_helpop.so
usr/lib/inspircd/modules/m_hidechans.so
usr/lib/inspircd/modules/m_hidelist.so
usr/lib/inspircd/modules/m_hidemode.so
usr/lib/inspircd/modules/m_hideoper.so
usr/lib/inspircd/modules/m_hostchange.so
usr/lib/inspircd/modules/m_hostcycle.so
usr/lib/inspircd/modules/m_httpd.so
usr/lib/inspircd/modules/m_httpd_acl.so
usr/lib/inspircd/modules/m_httpd_config.so
usr/lib/inspircd/modules/m_httpd_stats.so
usr/lib/inspircd/modules/m_ident.so
usr/lib/inspircd/modules/m_inviteexception.so
usr/lib/inspircd/modules/m_ircv3.so
usr/lib/inspircd/modules/m_ircv3_accounttag.so
usr/lib/inspircd/modules/m_ircv3_batch.so
usr/lib/inspircd/modules/m_ircv3_capnotify.so
usr/lib/inspircd/modules/m_ircv3_chghost.so
usr/lib/inspircd/modules/m_ircv3_ctctags.so
usr/lib/inspircd/modules/m_ircv3_echomessage.so
usr/lib/inspircd/modules/m_ircv3_invitenotify.so
usr/lib/inspircd/modules/m_ircv3_labeledresponse.so
usr/lib/inspircd/modules/m_ircv3_msgid.so
usr/lib/inspircd/modules/m_ircv3_servertime.so
usr/lib/inspircd/modules/m_ircv3_sts.so
usr/lib/inspircd/modules/m_joinflood.so
usr/lib/inspircd/modules/m_kicknorejoin.so
usr/lib/inspircd/modules/m_knock.so
usr/lib/inspircd/modules/m_ldap.so
usr/lib/inspircd/modules/m_ldapauth.so
usr/lib/inspircd/modules/m_ldapoper.so
usr/lib/inspircd/modules/m_lockserv.so
usr/lib/inspircd/modules/m_maphide.so
usr/lib/inspircd/modules/m_md5.so
usr/lib/inspircd/modules/m_messageflood.so
usr/lib/inspircd/modules/m_mlock.so
usr/lib/inspircd/modules/m_modenotice.so
usr/lib/inspircd/modules/m_monitor.so
usr/lib/inspircd/modules/m_muteban.so
usr/lib/inspircd/modules/m_mysql.so
usr/lib/inspircd/modules/m_namedmodes.so
usr/lib/inspircd/modules/m_namesx.so
usr/lib/inspircd/modules/m_nationalchars.so
usr/lib/inspircd/modules/m_nickflood.so
usr/lib/inspircd/modules/m_nicklock.so
usr/lib/inspircd/modules/m_noctcp.so
usr/lib/inspircd/modules/m_nokicks.so
usr/lib/inspircd/modules/m_nonicks.so
usr/lib/inspircd/modules/m_nonotice.so
usr/lib/inspircd/modules/m_nopartmsg.so
usr/lib/inspircd/modules/m_ojoin.so
usr/lib/inspircd/modules/m_operchans.so
usr/lib/inspircd/modules/m_operjoin.so
usr/lib/inspircd/modules/m_operlevels.so
usr/lib/inspircd/modules/m_operlog.so
usr/lib/inspircd/modules/m_opermodes.so
usr/lib/inspircd/modules/m_opermotd.so
usr/lib/inspircd/modules/m_operprefix.so
usr/lib/inspircd/modules/m_override.so
usr/lib/inspircd/modules/m_passforward.so
usr/lib/inspircd/modules/m_password_hash.so
usr/lib/inspircd/modules/m_pbkdf2.so
usr/lib/inspircd/modules/m_permchannels.so
usr/lib/inspircd/modules/m_pgsql.so
usr/lib/inspircd/modules/m_randquote.so
usr/lib/inspircd/modules/m_redirect.so
usr/lib/inspircd/modules/m_regex_glob.so
usr/lib/inspircd/modules/m_regex_pcre.so
usr/lib/inspircd/modules/m_regex_posix.so
usr/lib/inspircd/modules/m_regex_re2.so
usr/lib/inspircd/modules/m_regex_stdlib.so
usr/lib/inspircd/modules/m_regex_tre.so
usr/lib/inspircd/modules/m_remove.so
usr/lib/inspircd/modules/m_repeat.so
usr/lib/inspircd/modules/m_restrictchans.so
usr/lib/inspircd/modules/m_restrictmsg.so
usr/lib/inspircd/modules/m_rline.so
usr/lib/inspircd/modules/m_rmode.so
usr/lib/inspircd/modules/m_sajoin.so
usr/lib/inspircd/modules/m_sakick.so
usr/lib/inspircd/modules/m_samode.so
usr/lib/inspircd/modules/m_sanick.so
usr/lib/inspircd/modules/m_sapart.so
usr/lib/inspircd/modules/m_saquit.so
usr/lib/inspircd/modules/m_sasl.so
usr/lib/inspircd/modules/m_satopic.so
usr/lib/inspircd/modules/m_securelist.so
usr/lib/inspircd/modules/m_seenicks.so
usr/lib/inspircd/modules/m_serverban.so
usr/lib/inspircd/modules/m_services_account.so
usr/lib/inspircd/modules/m_servprotect.so
usr/lib/inspircd/modules/m_sethost.so
usr/lib/inspircd/modules/m_setident.so
usr/lib/inspircd/modules/m_setidle.so
usr/lib/inspircd/modules/m_setname.so
usr/lib/inspircd/modules/m_sha1.so
usr/lib/inspircd/modules/m_sha256.so
usr/lib/inspircd/modules/m_showfile.so
usr/lib/inspircd/modules/m_showwhois.so
usr/lib/inspircd/modules/m_shun.so
usr/lib/inspircd/modules/m_silence.so
usr/lib/inspircd/modules/m_spanningtree.so
usr/lib/inspircd/modules/m_sqlauth.so
usr/lib/inspircd/modules/m_sqlite3.so
usr/lib/inspircd/modules/m_sqloper.so
usr/lib/inspircd/modules/m_ssl_gnutls.so
usr/lib/inspircd/modules/m_ssl_mbedtls.so
usr/lib/inspircd/modules/m_ssl_openssl.so
usr/lib/inspircd/modules/m_sslinfo.so
usr/lib/inspircd/modules/m_sslmodes.so
usr/lib/inspircd/modules/m_sslrehashsignal.so
usr/lib/inspircd/modules/m_starttls.so
usr/lib/inspircd/modules/m_stripcolor.so
usr/lib/inspircd/modules/m_svshold.so
usr/lib/inspircd/modules/m_swhois.so
usr/lib/inspircd/modules/m_timedbans.so
usr/lib/inspircd/modules/m_tline.so
usr/lib/inspircd/modules/m_topiclock.so
usr/lib/inspircd/modules/m_uhnames.so
usr/lib/inspircd/modules/m_uninvite.so
usr/lib/inspircd/modules/m_userip.so
usr/lib/inspircd/modules/m_vhost.so
usr/lib/inspircd/modules/m_watch.so
usr/lib/inspircd/modules/m_websocket.so
usr/lib/inspircd/modules/m_xline_db.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/inspircd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/inspircd.conf
usr/share/
usr/share/inspircd/
usr/share/inspircd/examples/
usr/share/inspircd/examples/codepages/
usr/share/inspircd/examples/codepages/ascii.conf.example
usr/share/inspircd/examples/codepages/iso-8859-1.conf.example
usr/share/inspircd/examples/codepages/iso-8859-2.conf.example
usr/share/inspircd/examples/codepages/rfc1459.conf.example
usr/share/inspircd/examples/codepages/strict-rfc1459.conf.example
usr/share/inspircd/examples/filter.conf.example
usr/share/inspircd/examples/helpop.conf.example
usr/share/inspircd/examples/inspircd.conf.example
usr/share/inspircd/examples/links.conf.example
usr/share/inspircd/examples/modules.conf.example
usr/share/inspircd/examples/motd.txt.example
usr/share/inspircd/examples/opermotd.txt.example
usr/share/inspircd/examples/opers.conf.example
usr/share/inspircd/examples/providers/
usr/share/inspircd/examples/providers/dronebl.conf.example
usr/share/inspircd/examples/providers/efnet-rbl.conf.example
usr/share/inspircd/examples/providers/irccloud.conf.example
usr/share/inspircd/examples/providers/torexit.conf.example
usr/share/inspircd/examples/quotes.txt.example
usr/share/inspircd/examples/services/
usr/share/inspircd/examples/services/anope.conf.example
usr/share/inspircd/examples/services/atheme.conf.example
usr/share/inspircd/examples/services/generic.conf.example
usr/share/inspircd/examples/sql/
usr/share/inspircd/examples/sql/sqloper.mysql.sql
usr/share/inspircd/examples/sql/sqloper.pgsql.sql
usr/share/inspircd/examples/sql/sqloper.sqlite3.sql
var/
var/lib/
var/lib/inspircd/
var/log/
var/log/inspircd/
