# files/cockroachdb-bin-21.2.9-1/files
%FILES%
etc/
etc/default/
etc/default/cockroach
usr/
usr/bin/
usr/bin/cockroach
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/cockroach.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/cockroach.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/cockroach.conf
usr/local/
usr/local/lib/
usr/local/lib/cockroach/
usr/local/lib/cockroach/libgeos.so
usr/local/lib/cockroach/libgeos_c.so
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/cockroach
usr/share/licenses/
usr/share/licenses/cockroachdb-bin/
usr/share/licenses/cockroachdb-bin/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/cockroach-auth-session-list.1.gz
usr/share/man/man1/cockroach-auth-session-login.1.gz
usr/share/man/man1/cockroach-auth-session-logout.1.gz
usr/share/man/man1/cockroach-auth-session.1.gz
usr/share/man/man1/cockroach-cert-create-ca.1.gz
usr/share/man/man1/cockroach-cert-create-client-ca.1.gz
usr/share/man/man1/cockroach-cert-create-client.1.gz
usr/share/man/man1/cockroach-cert-create-node.1.gz
usr/share/man/man1/cockroach-cert-list.1.gz
usr/share/man/man1/cockroach-cert.1.gz
usr/share/man/man1/cockroach-connect-init.1.gz
usr/share/man/man1/cockroach-connect-join.1.gz
usr/share/man/man1/cockroach-connect.1.gz
usr/share/man/man1/cockroach-convert-url.1.gz
usr/share/man/man1/cockroach-debug-backup-export.1.gz
usr/share/man/man1/cockroach-debug-backup-list-backups.1.gz
usr/share/man/man1/cockroach-debug-backup-list-incremental.1.gz
usr/share/man/man1/cockroach-debug-backup-show.1.gz
usr/share/man/man1/cockroach-debug-backup.1.gz
usr/share/man/man1/cockroach-debug-ballast.1.gz
usr/share/man/man1/cockroach-debug-check-log-config.1.gz
usr/share/man/man1/cockroach-debug-check-store.1.gz
usr/share/man/man1/cockroach-debug-compact.1.gz
usr/share/man/man1/cockroach-debug-decode-key.1.gz
usr/share/man/man1/cockroach-debug-decode-proto.1.gz
usr/share/man/man1/cockroach-debug-decode-value.1.gz
usr/share/man/man1/cockroach-debug-doctor-examine-cluster.1.gz
usr/share/man/man1/cockroach-debug-doctor-examine-zipdir.1.gz
usr/share/man/man1/cockroach-debug-doctor-examine.1.gz
usr/share/man/man1/cockroach-debug-doctor-recreate-cluster.1.gz
usr/share/man/man1/cockroach-debug-doctor-recreate-zipdir.1.gz
usr/share/man/man1/cockroach-debug-doctor-recreate.1.gz
usr/share/man/man1/cockroach-debug-doctor.1.gz
usr/share/man/man1/cockroach-debug-encryption-active-key.1.gz
usr/share/man/man1/cockroach-debug-encryption-status.1.gz
usr/share/man/man1/cockroach-debug-env.1.gz
usr/share/man/man1/cockroach-debug-estimate-gc.1.gz
usr/share/man/man1/cockroach-debug-gossip-values.1.gz
usr/share/man/man1/cockroach-debug-intent-count.1.gz
usr/share/man/man1/cockroach-debug-job-trace.1.gz
usr/share/man/man1/cockroach-debug-keys.1.gz
usr/share/man/man1/cockroach-debug-list-files.1.gz
usr/share/man/man1/cockroach-debug-merge-logs.1.gz
usr/share/man/man1/cockroach-debug-pebble-db-check.1.gz
usr/share/man/man1/cockroach-debug-pebble-db-checkpoint.1.gz
usr/share/man/man1/cockroach-debug-pebble-db-get.1.gz
usr/share/man/man1/cockroach-debug-pebble-db-lsm.1.gz
usr/share/man/man1/cockroach-debug-pebble-db-properties.1.gz
usr/share/man/man1/cockroach-debug-pebble-db-scan.1.gz
usr/share/man/man1/cockroach-debug-pebble-db-set.1.gz
usr/share/man/man1/cockroach-debug-pebble-db-space.1.gz
usr/share/man/man1/cockroach-debug-pebble-db.1.gz
usr/share/man/man1/cockroach-debug-pebble-find.1.gz
usr/share/man/man1/cockroach-debug-pebble-lsm.1.gz
usr/share/man/man1/cockroach-debug-pebble-manifest-check.1.gz
usr/share/man/man1/cockroach-debug-pebble-manifest-dump.1.gz
usr/share/man/man1/cockroach-debug-pebble-manifest-summarize.1.gz
usr/share/man/man1/cockroach-debug-pebble-manifest.1.gz
usr/share/man/man1/cockroach-debug-pebble-sstable-check.1.gz
usr/share/man/man1/cockroach-debug-pebble-sstable-layout.1.gz
usr/share/man/man1/cockroach-debug-pebble-sstable-properties.1.gz
usr/share/man/man1/cockroach-debug-pebble-sstable-scan.1.gz
usr/share/man/man1/cockroach-debug-pebble-sstable-space.1.gz
usr/share/man/man1/cockroach-debug-pebble-sstable.1.gz
usr/share/man/man1/cockroach-debug-pebble-wal-dump.1.gz
usr/share/man/man1/cockroach-debug-pebble-wal.1.gz
usr/share/man/man1/cockroach-debug-pebble.1.gz
usr/share/man/man1/cockroach-debug-raft-log.1.gz
usr/share/man/man1/cockroach-debug-range-data.1.gz
usr/share/man/man1/cockroach-debug-range-descriptors.1.gz
usr/share/man/man1/cockroach-debug-reset-quorum.1.gz
usr/share/man/man1/cockroach-debug-statement-bundle-recreate.1.gz
usr/share/man/man1/cockroach-debug-statement-bundle.1.gz
usr/share/man/man1/cockroach-debug-synctest.1.gz
usr/share/man/man1/cockroach-debug-tsdump.1.gz
usr/share/man/man1/cockroach-debug-unsafe-remove-dead-replicas.1.gz
usr/share/man/man1/cockroach-debug-zip.1.gz
usr/share/man/man1/cockroach-debug.1.gz
usr/share/man/man1/cockroach-demo-bank.1.gz
usr/share/man/man1/cockroach-demo-intro.1.gz
usr/share/man/man1/cockroach-demo-kv.1.gz
usr/share/man/man1/cockroach-demo-movr.1.gz
usr/share/man/man1/cockroach-demo-startrek.1.gz
usr/share/man/man1/cockroach-demo-tpcc.1.gz
usr/share/man/man1/cockroach-demo-ycsb.1.gz
usr/share/man/man1/cockroach-demo.1.gz
usr/share/man/man1/cockroach-gen-autocomplete.1.gz
usr/share/man/man1/cockroach-gen-encryption-key.1.gz
usr/share/man/man1/cockroach-gen-example-data-bank.1.gz
usr/share/man/man1/cockroach-gen-example-data-bulkingest.1.gz
usr/share/man/man1/cockroach-gen-example-data-intro.1.gz
usr/share/man/man1/cockroach-gen-example-data-kv.1.gz
usr/share/man/man1/cockroach-gen-example-data-movr.1.gz
usr/share/man/man1/cockroach-gen-example-data-startrek.1.gz
usr/share/man/man1/cockroach-gen-example-data-tpcc.1.gz
usr/share/man/man1/cockroach-gen-example-data-tpch.1.gz
usr/share/man/man1/cockroach-gen-example-data-ycsb.1.gz
usr/share/man/man1/cockroach-gen-example-data.1.gz
usr/share/man/man1/cockroach-gen-haproxy.1.gz
usr/share/man/man1/cockroach-gen-man.1.gz
usr/share/man/man1/cockroach-gen-settings-list.1.gz
usr/share/man/man1/cockroach-gen.1.gz
usr/share/man/man1/cockroach-import-db.1.gz
usr/share/man/man1/cockroach-import-table.1.gz
usr/share/man/man1/cockroach-import.1.gz
usr/share/man/man1/cockroach-init.1.gz
usr/share/man/man1/cockroach-mt-cert-create-tenant-client-ca.1.gz
usr/share/man/man1/cockroach-mt-cert-create-tenant-client.1.gz
usr/share/man/man1/cockroach-node-decommission.1.gz
usr/share/man/man1/cockroach-node-drain.1.gz
usr/share/man/man1/cockroach-node-ls.1.gz
usr/share/man/man1/cockroach-node-recommission.1.gz
usr/share/man/man1/cockroach-node-status.1.gz
usr/share/man/man1/cockroach-node.1.gz
usr/share/man/man1/cockroach-nodelocal-upload.1.gz
usr/share/man/man1/cockroach-nodelocal.1.gz
usr/share/man/man1/cockroach-sql.1.gz
usr/share/man/man1/cockroach-sqlfmt.1.gz
usr/share/man/man1/cockroach-start-single-node.1.gz
usr/share/man/man1/cockroach-start.1.gz
usr/share/man/man1/cockroach-statement-diag-cancel.1.gz
usr/share/man/man1/cockroach-statement-diag-delete.1.gz
usr/share/man/man1/cockroach-statement-diag-download.1.gz
usr/share/man/man1/cockroach-statement-diag-list.1.gz
usr/share/man/man1/cockroach-statement-diag.1.gz
usr/share/man/man1/cockroach-userfile-delete.1.gz
usr/share/man/man1/cockroach-userfile-get.1.gz
usr/share/man/man1/cockroach-userfile-list.1.gz
usr/share/man/man1/cockroach-userfile-upload.1.gz
usr/share/man/man1/cockroach-userfile.1.gz
usr/share/man/man1/cockroach-version.1.gz
usr/share/man/man1/cockroach-workload-init-bank.1.gz
usr/share/man/man1/cockroach-workload-init-intro.1.gz
usr/share/man/man1/cockroach-workload-init-kv.1.gz
usr/share/man/man1/cockroach-workload-init-movr.1.gz
usr/share/man/man1/cockroach-workload-init-startrek.1.gz
usr/share/man/man1/cockroach-workload-init-tpcc.1.gz
usr/share/man/man1/cockroach-workload-init-ycsb.1.gz
usr/share/man/man1/cockroach-workload-init.1.gz
usr/share/man/man1/cockroach-workload-run-bank.1.gz
usr/share/man/man1/cockroach-workload-run-kv.1.gz
usr/share/man/man1/cockroach-workload-run-movr.1.gz
usr/share/man/man1/cockroach-workload-run-tpcc.1.gz
usr/share/man/man1/cockroach-workload-run-ycsb.1.gz
usr/share/man/man1/cockroach-workload-run.1.gz
usr/share/man/man1/cockroach-workload.1.gz
usr/share/man/man1/cockroach.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_cockroach
