# files/betterlockscreen-4.0.4-1/files
%FILES%
usr/
usr/bin/
usr/bin/betterlockscreen
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/betterlockscreen@.service
usr/share/
usr/share/doc/
usr/share/doc/betterlockscreen/
usr/share/doc/betterlockscreen/examples/
usr/share/doc/betterlockscreen/examples/betterlockscreenrc
usr/share/licenses/
usr/share/licenses/betterlockscreen/
usr/share/licenses/betterlockscreen/LICENSE
