# files/popura-git-v0.4.3+popura1.r2.gbbd162c-1/files
%FILES%
usr/
usr/bin/
usr/bin/yggdrasil
usr/bin/yggdrasilctl
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/yggdrasil-default-config.service
usr/lib/systemd/system/yggdrasil.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/popura.conf
usr/share/
usr/share/licenses/
usr/share/licenses/popura/
usr/share/licenses/popura/LICENSE
