# files/kernel-modules-hook-0.1.7-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/linux-modules-cleanup.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/linux-modules-cleanup.conf
usr/share/
usr/share/libalpm/
usr/share/libalpm/hooks/
usr/share/libalpm/hooks/10-linux-modules-post.hook
usr/share/libalpm/hooks/10-linux-modules-pre.hook
usr/share/licenses/
usr/share/licenses/kernel-modules-hook/
usr/share/licenses/kernel-modules-hook/UNLICENSE
