# files/input-remapper-git-1.4.0.r7.g6029aac-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/inputremapper.Control.conf
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/input-remapper-autoload.desktop
usr/
usr/bin/
usr/bin/input-remapper-control
usr/bin/input-remapper-gtk
usr/bin/input-remapper-helper
usr/bin/input-remapper-service
usr/bin/key-mapper-control
usr/bin/key-mapper-gtk
usr/bin/key-mapper-service
usr/lib/
usr/lib/python3.10/
usr/lib/python3.10/site-packages/
usr/lib/python3.10/site-packages/input_remapper-1.4.0-py3.10.egg-info/
usr/lib/python3.10/site-packages/input_remapper-1.4.0-py3.10.egg-info/PKG-INFO
usr/lib/python3.10/site-packages/input_remapper-1.4.0-py3.10.egg-info/SOURCES.txt
usr/lib/python3.10/site-packages/input_remapper-1.4.0-py3.10.egg-info/dependency_links.txt
usr/lib/python3.10/site-packages/input_remapper-1.4.0-py3.10.egg-info/requires.txt
usr/lib/python3.10/site-packages/input_remapper-1.4.0-py3.10.egg-info/top_level.txt
usr/lib/python3.10/site-packages/inputremapper/
usr/lib/python3.10/site-packages/inputremapper/__init__.py
usr/lib/python3.10/site-packages/inputremapper/__pycache__/
usr/lib/python3.10/site-packages/inputremapper/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/__pycache__/commit_hash.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/__pycache__/daemon.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/__pycache__/event_combination.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/__pycache__/exceptions.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/__pycache__/groups.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/__pycache__/input_event.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/__pycache__/logger.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/__pycache__/user.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/__pycache__/utils.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/commit_hash.py
usr/lib/python3.10/site-packages/inputremapper/configs/
usr/lib/python3.10/site-packages/inputremapper/configs/__init__.py
usr/lib/python3.10/site-packages/inputremapper/configs/__pycache__/
usr/lib/python3.10/site-packages/inputremapper/configs/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/configs/__pycache__/base_config.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/configs/__pycache__/data.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/configs/__pycache__/global_config.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/configs/__pycache__/migrations.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/configs/__pycache__/paths.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/configs/__pycache__/preset.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/configs/__pycache__/system_mapping.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/configs/base_config.py
usr/lib/python3.10/site-packages/inputremapper/configs/data.py
usr/lib/python3.10/site-packages/inputremapper/configs/global_config.py
usr/lib/python3.10/site-packages/inputremapper/configs/migrations.py
usr/lib/python3.10/site-packages/inputremapper/configs/paths.py
usr/lib/python3.10/site-packages/inputremapper/configs/preset.py
usr/lib/python3.10/site-packages/inputremapper/configs/system_mapping.py
usr/lib/python3.10/site-packages/inputremapper/daemon.py
usr/lib/python3.10/site-packages/inputremapper/event_combination.py
usr/lib/python3.10/site-packages/inputremapper/exceptions.py
usr/lib/python3.10/site-packages/inputremapper/groups.py
usr/lib/python3.10/site-packages/inputremapper/gui/
usr/lib/python3.10/site-packages/inputremapper/gui/__init__.py
usr/lib/python3.10/site-packages/inputremapper/gui/__pycache__/
usr/lib/python3.10/site-packages/inputremapper/gui/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/gui/__pycache__/active_preset.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/gui/__pycache__/helper.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/gui/__pycache__/reader.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/gui/__pycache__/user_interface.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/gui/__pycache__/utils.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/gui/active_preset.py
usr/lib/python3.10/site-packages/inputremapper/gui/editor/
usr/lib/python3.10/site-packages/inputremapper/gui/editor/__init__.py
usr/lib/python3.10/site-packages/inputremapper/gui/editor/__pycache__/
usr/lib/python3.10/site-packages/inputremapper/gui/editor/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/gui/editor/__pycache__/autocompletion.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/gui/editor/__pycache__/editor.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/gui/editor/autocompletion.py
usr/lib/python3.10/site-packages/inputremapper/gui/editor/editor.py
usr/lib/python3.10/site-packages/inputremapper/gui/helper.py
usr/lib/python3.10/site-packages/inputremapper/gui/reader.py
usr/lib/python3.10/site-packages/inputremapper/gui/user_interface.py
usr/lib/python3.10/site-packages/inputremapper/gui/utils.py
usr/lib/python3.10/site-packages/inputremapper/injection/
usr/lib/python3.10/site-packages/inputremapper/injection/__init__.py
usr/lib/python3.10/site-packages/inputremapper/injection/__pycache__/
usr/lib/python3.10/site-packages/inputremapper/injection/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/__pycache__/consumer_control.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/__pycache__/context.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/__pycache__/global_uinputs.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/__pycache__/injector.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/__pycache__/numlock.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/consumer_control.py
usr/lib/python3.10/site-packages/inputremapper/injection/consumers/
usr/lib/python3.10/site-packages/inputremapper/injection/consumers/__init__.py
usr/lib/python3.10/site-packages/inputremapper/injection/consumers/__pycache__/
usr/lib/python3.10/site-packages/inputremapper/injection/consumers/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/consumers/__pycache__/consumer.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/consumers/__pycache__/joystick_to_mouse.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/consumers/__pycache__/keycode_mapper.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/consumers/consumer.py
usr/lib/python3.10/site-packages/inputremapper/injection/consumers/joystick_to_mouse.py
usr/lib/python3.10/site-packages/inputremapper/injection/consumers/keycode_mapper.py
usr/lib/python3.10/site-packages/inputremapper/injection/context.py
usr/lib/python3.10/site-packages/inputremapper/injection/global_uinputs.py
usr/lib/python3.10/site-packages/inputremapper/injection/injector.py
usr/lib/python3.10/site-packages/inputremapper/injection/macros/
usr/lib/python3.10/site-packages/inputremapper/injection/macros/__init__.py
usr/lib/python3.10/site-packages/inputremapper/injection/macros/__pycache__/
usr/lib/python3.10/site-packages/inputremapper/injection/macros/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/macros/__pycache__/macro.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/macros/__pycache__/parse.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/injection/macros/macro.py
usr/lib/python3.10/site-packages/inputremapper/injection/macros/parse.py
usr/lib/python3.10/site-packages/inputremapper/injection/numlock.py
usr/lib/python3.10/site-packages/inputremapper/input_event.py
usr/lib/python3.10/site-packages/inputremapper/ipc/
usr/lib/python3.10/site-packages/inputremapper/ipc/__init__.py
usr/lib/python3.10/site-packages/inputremapper/ipc/__pycache__/
usr/lib/python3.10/site-packages/inputremapper/ipc/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/ipc/__pycache__/pipe.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/ipc/__pycache__/shared_dict.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/ipc/__pycache__/socket.cpython-310.pyc
usr/lib/python3.10/site-packages/inputremapper/ipc/pipe.py
usr/lib/python3.10/site-packages/inputremapper/ipc/shared_dict.py
usr/lib/python3.10/site-packages/inputremapper/ipc/socket.py
usr/lib/python3.10/site-packages/inputremapper/logger.py
usr/lib/python3.10/site-packages/inputremapper/user.py
usr/lib/python3.10/site-packages/inputremapper/utils.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/input-remapper.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/99-input-remapper.rules
usr/share/
usr/share/applications/
usr/share/applications/input-remapper.desktop
usr/share/input-remapper/
usr/share/input-remapper/99-input-remapper.rules
usr/share/input-remapper/input-remapper-autoload.desktop
usr/share/input-remapper/input-remapper-large.png
usr/share/input-remapper/input-remapper.desktop
usr/share/input-remapper/input-remapper.glade
usr/share/input-remapper/input-remapper.policy
usr/share/input-remapper/input-remapper.service
usr/share/input-remapper/input-remapper.svg
usr/share/input-remapper/inputremapper.Control.conf
usr/share/input-remapper/lang/
usr/share/input-remapper/lang/it_IT/
usr/share/input-remapper/lang/it_IT/LC_MESSAGES/
usr/share/input-remapper/lang/it_IT/LC_MESSAGES/input-remapper.mo
usr/share/input-remapper/lang/sk_SK/
usr/share/input-remapper/lang/sk_SK/LC_MESSAGES/
usr/share/input-remapper/lang/sk_SK/LC_MESSAGES/input-remapper.mo
usr/share/input-remapper/style.css
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/input-remapper.policy
