# files/timeshift-22.06.1-1/files
%FILES%
etc/
etc/systemd/
etc/systemd/system/
etc/systemd/system/grub-btrfs.path
etc/timeshift/
etc/timeshift/default.json
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/snapshot-detect.desktop
usr/
usr/bin/
usr/bin/snapshot-detect
usr/bin/timeshift
usr/bin/timeshift-gtk
usr/bin/timeshift-launcher
usr/share/
usr/share/appdata/
usr/share/applications/
usr/share/applications/timeshift-gtk.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/apps/
usr/share/icons/hicolor/128x128/apps/timeshift.png
usr/share/icons/hicolor/16x16/
usr/share/icons/hicolor/16x16/apps/
usr/share/icons/hicolor/16x16/apps/timeshift.png
usr/share/icons/hicolor/22x22/
usr/share/icons/hicolor/22x22/apps/
usr/share/icons/hicolor/22x22/apps/timeshift.png
usr/share/icons/hicolor/24x24/
usr/share/icons/hicolor/24x24/apps/
usr/share/icons/hicolor/24x24/apps/timeshift.png
usr/share/icons/hicolor/32x32/
usr/share/icons/hicolor/32x32/apps/
usr/share/icons/hicolor/32x32/apps/timeshift.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/timeshift.png
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/apps/
usr/share/icons/hicolor/64x64/apps/timeshift.png
usr/share/icons/hicolor/96x96/
usr/share/icons/hicolor/96x96/apps/
usr/share/icons/hicolor/96x96/apps/timeshift.png
usr/share/locale/
usr/share/locale/am/
usr/share/locale/am/LC_MESSAGES/
usr/share/locale/am/LC_MESSAGES/timeshift.mo
usr/share/locale/ar/
usr/share/locale/ar/LC_MESSAGES/
usr/share/locale/ar/LC_MESSAGES/timeshift.mo
usr/share/locale/az/
usr/share/locale/az/LC_MESSAGES/
usr/share/locale/az/LC_MESSAGES/timeshift.mo
usr/share/locale/bg/
usr/share/locale/bg/LC_MESSAGES/
usr/share/locale/bg/LC_MESSAGES/timeshift.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/timeshift.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/timeshift.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/timeshift.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/timeshift.mo
usr/share/locale/el/
usr/share/locale/el/LC_MESSAGES/
usr/share/locale/el/LC_MESSAGES/timeshift.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/timeshift.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/timeshift.mo
usr/share/locale/et/
usr/share/locale/et/LC_MESSAGES/
usr/share/locale/et/LC_MESSAGES/timeshift.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/timeshift.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/timeshift.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/timeshift.mo
usr/share/locale/he/
usr/share/locale/he/LC_MESSAGES/
usr/share/locale/he/LC_MESSAGES/timeshift.mo
usr/share/locale/hi/
usr/share/locale/hi/LC_MESSAGES/
usr/share/locale/hi/LC_MESSAGES/timeshift.mo
usr/share/locale/hr/
usr/share/locale/hr/LC_MESSAGES/
usr/share/locale/hr/LC_MESSAGES/timeshift.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/timeshift.mo
usr/share/locale/ia/
usr/share/locale/ia/LC_MESSAGES/
usr/share/locale/ia/LC_MESSAGES/timeshift.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/timeshift.mo
usr/share/locale/is/
usr/share/locale/is/LC_MESSAGES/
usr/share/locale/is/LC_MESSAGES/timeshift.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/timeshift.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/timeshift.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/timeshift.mo
usr/share/locale/lt/
usr/share/locale/lt/LC_MESSAGES/
usr/share/locale/lt/LC_MESSAGES/timeshift.mo
usr/share/locale/nb/
usr/share/locale/nb/LC_MESSAGES/
usr/share/locale/nb/LC_MESSAGES/timeshift.mo
usr/share/locale/ne/
usr/share/locale/ne/LC_MESSAGES/
usr/share/locale/ne/LC_MESSAGES/timeshift.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/timeshift.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/timeshift.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/timeshift.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/timeshift.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/timeshift.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/timeshift.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/timeshift.mo
usr/share/locale/sr/
usr/share/locale/sr/LC_MESSAGES/
usr/share/locale/sr/LC_MESSAGES/timeshift.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/timeshift.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/timeshift.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/timeshift.mo
usr/share/locale/vi/
usr/share/locale/vi/LC_MESSAGES/
usr/share/locale/vi/LC_MESSAGES/timeshift.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/timeshift.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/timeshift.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/timeshift.1.gz
usr/share/metainfo/
usr/share/metainfo/timeshift.appdata.xml
usr/share/pixmaps/
usr/share/pixmaps/timeshift.png
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/in.teejeetech.pkexec.timeshift.policy
usr/share/timeshift/
usr/share/timeshift/images/
usr/share/timeshift/images/clock.png
usr/share/timeshift/images/disk.png
usr/share/timeshift/images/document-open-recent-symbolic.svg
usr/share/timeshift/images/document-save-symbolic.svg
usr/share/timeshift/images/donate.png
usr/share/timeshift/images/drive-harddisk.svg
usr/share/timeshift/images/edit-delete-symbolic.svg
usr/share/timeshift/images/edit-delete.png
usr/share/timeshift/images/emblem-default-symbolic.svg
usr/share/timeshift/images/exclude.png
usr/share/timeshift/images/folder-symbolic.svg
usr/share/timeshift/images/gtk-file.png
usr/share/timeshift/images/include.png
usr/share/timeshift/images/item-blue.png
usr/share/timeshift/images/item-gray.png
usr/share/timeshift/images/item-green.png
usr/share/timeshift/images/item-red.png
usr/share/timeshift/images/item-yellow.png
usr/share/timeshift/images/list-add.png
usr/share/timeshift/images/list-remove.png
usr/share/timeshift/images/locked.png
usr/share/timeshift/images/media-optical.png
usr/share/timeshift/images/open-menu-symbolic.svg
usr/share/timeshift/images/preferences-system-symbolic.svg
usr/share/timeshift/images/timeshift-shield-high.svg
usr/share/timeshift/images/timeshift-shield-low.svg
usr/share/timeshift/images/timeshift-shield-med.svg
usr/share/timeshift/images/unlocked.png
