# files/tor-git-0.4.7.0.alpha.dev.r96.ge247aab-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/tor
etc/tor/
etc/tor/torrc
etc/tor/torrc.d/
etc/tor/torrc.d/bridge
etc/tor/torrc.d/isolation
etc/tor/torrc.d/nodes
etc/tor/torrc.d/transparent_proxy
etc/tor/torrc.sample
usr/
usr/bin/
usr/bin/tor
usr/bin/tor-gencert
usr/bin/tor-print-ed-signing-cert
usr/bin/tor-resolve
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tor.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/tor.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/tor.conf
usr/share/
usr/share/licenses/
usr/share/licenses/tor-git/
usr/share/licenses/tor-git/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/tor-gencert.1.gz
usr/share/man/man1/tor-print-ed-signing-cert.1.gz
usr/share/man/man1/tor-resolve.1.gz
usr/share/man/man1/tor.1.gz
usr/share/man/man1/torify.1.gz
usr/share/tor/
usr/share/tor/geoip
usr/share/tor/geoip6
