# files/keymapper-git-r173.8d59442-1/files
%FILES%
usr/
usr/bin/
usr/bin/keymapper
usr/bin/keymapperd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/keymapperd.service
usr/share/
usr/share/applications/
usr/share/applications/keymapper.desktop
usr/share/gnome-shell/
usr/share/gnome-shell/extensions/
usr/share/gnome-shell/extensions/keymapper@houmain.github.com/
usr/share/gnome-shell/extensions/keymapper@houmain.github.com/extension.js
usr/share/gnome-shell/extensions/keymapper@houmain.github.com/metadata.json
