# files/cfs-zen-tweaks-1.2.0-1/files
%FILES%
usr/
usr/lib/
usr/lib/cfs-zen-tweaks/
usr/lib/cfs-zen-tweaks/set-cfs-zen-tweaks.bash
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/set-cfs-tweaks.service
