# files/anbox-git-1:r1358.84f0268-7/files
%FILES%
etc/
etc/xdg/
etc/xdg/menus/
etc/xdg/menus/applications-merged/
etc/xdg/menus/applications-merged/anbox-android.menu
usr/
usr/bin/
usr/bin/anbox
usr/include/
usr/include/backward.hpp
usr/lib/
usr/lib/backward/
usr/lib/backward/BackwardConfig.cmake
usr/lib/systemd/
usr/lib/systemd/network/
usr/lib/systemd/network/80-anbox-bridge.netdev
usr/lib/systemd/network/80-anbox-bridge.network
usr/lib/systemd/system/
usr/lib/systemd/system/anbox-container-manager.service
usr/lib/systemd/system/dev-binderfs.mount
usr/lib/systemd/user/
usr/lib/systemd/user/anbox-session-manager.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/99-anbox.rules
usr/share/
usr/share/anbox/
usr/share/anbox/ui/
usr/share/anbox/ui/loading-screen.png
usr/share/applications/
usr/share/applications/anbox.desktop
usr/share/dbus-1/
usr/share/dbus-1/services/
usr/share/dbus-1/services/org.anbox.service
usr/share/desktop-directories/
usr/share/desktop-directories/anbox-android.directory
usr/share/pixmaps/
usr/share/pixmaps/anbox.png
