# files/btrbk-0.32.4-1/files
%FILES%
etc/
etc/btrbk/
etc/btrbk/btrbk.conf.example
usr/
usr/bin/
usr/bin/btrbk
usr/bin/lsbtr
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/btrbk.service
usr/lib/systemd/system/btrbk.timer
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/btrbk
usr/share/bash-completion/completions/lsbtr
usr/share/btrbk/
usr/share/btrbk/scripts/
usr/share/btrbk/scripts/btrbk-mail
usr/share/btrbk/scripts/btrbk-verify
usr/share/btrbk/scripts/btrbk_restore_raw.py
usr/share/btrbk/scripts/kdf_pbkdf2.py
usr/share/btrbk/scripts/raw_suffix2sidecar
usr/share/btrbk/scripts/ssh_filter_btrbk.sh
usr/share/doc/
usr/share/doc/btrbk/
usr/share/doc/btrbk/ChangeLog
usr/share/doc/btrbk/FAQ.md
usr/share/doc/btrbk/README.md
usr/share/doc/btrbk/upgrade_to_v0.23.0.md
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/btrbk.1.gz
usr/share/man/man1/lsbtr.1.gz
usr/share/man/man1/ssh_filter_btrbk.1.gz
usr/share/man/man5/
usr/share/man/man5/btrbk.conf.5.gz
