# files/freeipa-client-4.10.0-1/files
%FILES%
etc/
etc/bash_completion.d/
etc/bash_completion.d/ipa
etc/krb5.conf.d/
usr/
usr/bin/
usr/bin/ipa
usr/bin/ipa-certupdate
usr/bin/ipa-client-automount
usr/bin/ipa-client-install
usr/bin/ipa-getkeytab
usr/bin/ipa-join
usr/bin/ipa-rmkeytab
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nis-domainname.service
usr/share/
usr/share/doc/
usr/share/doc/freeipa-client/
usr/share/doc/freeipa-client/Contributors.txt
usr/share/doc/freeipa-client/README.md
usr/share/libalpm/
usr/share/libalpm/hooks/
usr/share/libalpm/hooks/freeipa-client-update-sshd_config.hook
usr/share/libalpm/scripts/
usr/share/libalpm/scripts/freeipa-client-update-sshd_config
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/ipa-certupdate.1.gz
usr/share/man/man1/ipa-client-automount.1.gz
usr/share/man/man1/ipa-client-install.1.gz
usr/share/man/man1/ipa-getkeytab.1.gz
usr/share/man/man1/ipa-join.1.gz
usr/share/man/man1/ipa-rmkeytab.1.gz
usr/share/man/man1/ipa.1.gz
