# files/clamav-git-0.99.gd57d89865-1/files
%FILES%
etc/
etc/clamav/
etc/clamav/clamav-milter.conf.sample
etc/clamav/clamd.conf
etc/clamav/clamd.conf.sample
etc/clamav/freshclam.conf
etc/clamav/freshclam.conf.sample
etc/logrotate.d/
etc/logrotate.d/clamav
usr/
usr/bin/
usr/bin/clamav-config
usr/bin/clamav-milter
usr/bin/clambc
usr/bin/clamconf
usr/bin/clamd
usr/bin/clamdscan
usr/bin/clamdtop
usr/bin/clamonacc
usr/bin/clamscan
usr/bin/clamsubmit
usr/bin/freshclam
usr/bin/sigtool
usr/include/
usr/include/clamav-types.h
usr/include/clamav-version.h
usr/include/clamav.h
usr/include/libfreshclam.h
usr/lib/
usr/lib/libclamav.so
usr/lib/libclamav.so.9
usr/lib/libclamav.so.9.0.6
usr/lib/libclammspack.so
usr/lib/libclammspack.so.0
usr/lib/libclammspack.so.0.1.0
usr/lib/libclamunrar.so
usr/lib/libclamunrar.so.9
usr/lib/libclamunrar.so.9.0.6
usr/lib/libclamunrar_iface.so
usr/lib/libclamunrar_iface.so.9
usr/lib/libclamunrar_iface.so.9.0.6
usr/lib/libfreshclam.so
usr/lib/libfreshclam.so.2
usr/lib/libfreshclam.so.2.0.2
usr/lib/pkgconfig/
usr/lib/pkgconfig/libclamav.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/clamav-clamonacc.service
usr/lib/systemd/system/clamav-daemon.service
usr/lib/systemd/system/clamav-daemon.socket
usr/lib/systemd/system/clamav-freshclam.service
usr/lib/systemd/system/clamd.service
usr/lib/systemd/system/freshclamd.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/clamav.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/clambc.1.gz
usr/share/man/man1/clamconf.1.gz
usr/share/man/man1/clamdscan.1.gz
usr/share/man/man1/clamdtop.1.gz
usr/share/man/man1/clamscan.1.gz
usr/share/man/man1/clamsubmit.1.gz
usr/share/man/man1/freshclam.1.gz
usr/share/man/man1/sigtool.1.gz
usr/share/man/man5/
usr/share/man/man5/clamav-milter.conf.5.gz
usr/share/man/man5/clamd.conf.5.gz
usr/share/man/man5/freshclam.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/clamav-milter.8.gz
usr/share/man/man8/clamd.8.gz
usr/share/man/man8/clamonacc.8.gz
var/
var/lib/
var/lib/clamav/
var/log/
var/log/clamav/
