# files/xboxdrv-0.8.8-7/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/org.seul.Xboxdrv.conf
etc/default/
etc/default/xboxdrv
usr/
usr/bin/
usr/bin/xboxdrv
usr/bin/xboxdrvctl
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/xboxdrv.service
usr/share/
usr/share/doc/
usr/share/doc/xboxdrv/
usr/share/doc/xboxdrv/NEWS
usr/share/doc/xboxdrv/PROTOCOL
usr/share/doc/xboxdrv/README.md
usr/share/doc/xboxdrv/examples/
usr/share/doc/xboxdrv/examples/absolutemouse.xboxdrv
usr/share/doc/xboxdrv/examples/calibrate.macro
usr/share/doc/xboxdrv/examples/constfilter.xboxdrv
usr/share/doc/xboxdrv/examples/cyclekey.xboxdrv
usr/share/doc/xboxdrv/examples/default.xboxdrv
usr/share/doc/xboxdrv/examples/dosbox.xboxdrv
usr/share/doc/xboxdrv/examples/evdev.xboxdrv
usr/share/doc/xboxdrv/examples/fceux.xboxdrv
usr/share/doc/xboxdrv/examples/fceux2.xboxdrv
usr/share/doc/xboxdrv/examples/helloworld.macro
usr/share/doc/xboxdrv/examples/microsoft-ergonomic-4000.xboxdrv
usr/share/doc/xboxdrv/examples/mouse.xboxdrv
usr/share/doc/xboxdrv/examples/multi-controller.sh
usr/share/doc/xboxdrv/examples/playstation3.xboxdrv
usr/share/doc/xboxdrv/examples/ps4_usb.xboxdrv
usr/share/doc/xboxdrv/examples/responsecurve-generator.py
usr/share/doc/xboxdrv/examples/rotate.xboxdrv
usr/share/doc/xboxdrv/examples/sensitivity.xboxdrv
usr/share/doc/xboxdrv/examples/shift.xboxdrv
usr/share/doc/xboxdrv/examples/tiefighter.xboxdrv
usr/share/doc/xboxdrv/examples/two-sticks.xboxdrv
usr/share/doc/xboxdrv/examples/warsow.xboxdrv
usr/share/doc/xboxdrv/examples/xbmc.xboxdrv
usr/share/doc/xboxdrv/examples/xboxdrv-dualstrike
usr/share/doc/xboxdrv/examples/xboxdrv-dualstrike-mouse
usr/share/doc/xboxdrv/examples/xboxdrv-ps3-usb
usr/share/doc/xboxdrv/examples/xpad-wireless.xboxdrv
usr/share/doc/xboxdrv/examples/zsnes.xboxdrv
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/xboxdrv.1.gz
