# files/ums-10.8.2-1/files
%FILES%
opt/
opt/ums/
opt/ums/BUILD.md
opt/ums/CHANGELOG.txt
opt/ums/DummyInput.ass
opt/ums/DummyInput.jpg
opt/ums/INSTALL.txt
opt/ums/LICENSE.txt
opt/ums/README.md
opt/ums/UMS.conf
opt/ums/UMS.sh
opt/ums/VirtualFolders.conf
opt/ums/WEB.conf
opt/ums/database/
opt/ums/debug.log
opt/ums/documentation/
opt/ums/documentation/applications.html
opt/ums/documentation/avisynth.html
opt/ums/documentation/commandline.html
opt/ums/documentation/coreavc.html
opt/ums/documentation/css/
opt/ums/documentation/css/style.css
opt/ums/documentation/general_configuration.html
opt/ums/documentation/images/
opt/ums/documentation/images/autonegotiation.jpg
opt/ums/documentation/images/btn_donateCC_LG.gif
opt/ums/documentation/images/coreavc1.jpg
opt/ums/documentation/images/coreavc2.jpg
opt/ums/documentation/images/coreavc3.jpg
opt/ums/documentation/images/netzwerkverbindungen.jpg
opt/ums/documentation/images/pms-eng-general.jpg
opt/ums/documentation/images/pms-eng-general1.jpg
opt/ums/documentation/images/pms-eng-status1.jpg
opt/ums/documentation/images/pms-eng-status2.jpg
opt/ums/documentation/images/pms-engines1.jpg
opt/ums/documentation/images/pms-engines2.jpg
opt/ums/documentation/images/pms-searching.jpg
opt/ums/documentation/images/ps3-browsing.jpg
opt/ums/documentation/images/ps3-server.jpg
opt/ums/documentation/images/ps3screen.jpg
opt/ums/documentation/images/setupwin1.jpg
opt/ums/documentation/images/setupwin2.jpg
opt/ums/documentation/images/setupwin3.jpg
opt/ums/documentation/images/setupwin4.jpg
opt/ums/documentation/images/ssdp.jpg
opt/ums/documentation/images/topology_ps3.jpg
opt/ums/documentation/index.html
opt/ums/documentation/installation.html
opt/ums/documentation/introduction.html
opt/ums/documentation/links.html
opt/ums/documentation/navigation_share.html
opt/ums/documentation/networking.html
opt/ums/documentation/plugins.html
opt/ums/documentation/transcoding.html
opt/ums/documentation/videolan.html
opt/ums/jre8/
opt/ums/jre8/LICENSE
opt/ums/jre8/bin/
opt/ums/jre8/bin/java
opt/ums/jre8/bin/jjs
opt/ums/jre8/bin/keytool
opt/ums/jre8/bin/orbd
opt/ums/jre8/bin/pack200
opt/ums/jre8/bin/policytool
opt/ums/jre8/bin/rmid
opt/ums/jre8/bin/rmiregistry
opt/ums/jre8/bin/servertool
opt/ums/jre8/bin/tnameserv
opt/ums/jre8/bin/unpack200
opt/ums/jre8/lib/
opt/ums/jre8/lib/amd64/
opt/ums/jre8/lib/amd64/jli/
opt/ums/jre8/lib/amd64/jli/libjli.so
opt/ums/jre8/lib/amd64/jvm.cfg
opt/ums/jre8/lib/amd64/libawt.so
opt/ums/jre8/lib/amd64/libawt_headless.so
opt/ums/jre8/lib/amd64/libawt_xawt.so
opt/ums/jre8/lib/amd64/libdt_socket.so
opt/ums/jre8/lib/amd64/libfontmanager.so
opt/ums/jre8/lib/amd64/libhprof.so
opt/ums/jre8/lib/amd64/libinstrument.so
opt/ums/jre8/lib/amd64/libj2gss.so
opt/ums/jre8/lib/amd64/libj2pcsc.so
opt/ums/jre8/lib/amd64/libj2pkcs11.so
opt/ums/jre8/lib/amd64/libjaas_unix.so
opt/ums/jre8/lib/amd64/libjava.so
opt/ums/jre8/lib/amd64/libjava_crw_demo.so
opt/ums/jre8/lib/amd64/libjawt.so
opt/ums/jre8/lib/amd64/libjdwp.so
opt/ums/jre8/lib/amd64/libjpeg.so
opt/ums/jre8/lib/amd64/libjsdt.so
opt/ums/jre8/lib/amd64/libjsig.so
opt/ums/jre8/lib/amd64/libjsound.so
opt/ums/jre8/lib/amd64/libjsoundalsa.so
opt/ums/jre8/lib/amd64/liblcms.so
opt/ums/jre8/lib/amd64/libmanagement.so
opt/ums/jre8/lib/amd64/libmlib_image.so
opt/ums/jre8/lib/amd64/libnet.so
opt/ums/jre8/lib/amd64/libnio.so
opt/ums/jre8/lib/amd64/libnpt.so
opt/ums/jre8/lib/amd64/libsctp.so
opt/ums/jre8/lib/amd64/libsplashscreen.so
opt/ums/jre8/lib/amd64/libsunec.so
opt/ums/jre8/lib/amd64/libunpack.so
opt/ums/jre8/lib/amd64/libverify.so
opt/ums/jre8/lib/amd64/libzip.so
opt/ums/jre8/lib/amd64/server/
opt/ums/jre8/lib/amd64/server/Xusage.txt
opt/ums/jre8/lib/amd64/server/libjsig.so
opt/ums/jre8/lib/amd64/server/libjvm.so
opt/ums/jre8/lib/applet/
opt/ums/jre8/lib/calendars.properties
opt/ums/jre8/lib/charsets.jar
opt/ums/jre8/lib/classlist
opt/ums/jre8/lib/cmm/
opt/ums/jre8/lib/cmm/CIEXYZ.pf
opt/ums/jre8/lib/cmm/GRAY.pf
opt/ums/jre8/lib/cmm/LINEAR_RGB.pf
opt/ums/jre8/lib/cmm/PYCC.pf
opt/ums/jre8/lib/cmm/sRGB.pf
opt/ums/jre8/lib/content-types.properties
opt/ums/jre8/lib/currency.data
opt/ums/jre8/lib/ext/
opt/ums/jre8/lib/ext/cldrdata.jar
opt/ums/jre8/lib/ext/dnsns.jar
opt/ums/jre8/lib/ext/jaccess.jar
opt/ums/jre8/lib/ext/localedata.jar
opt/ums/jre8/lib/ext/meta-index
opt/ums/jre8/lib/ext/nashorn.jar
opt/ums/jre8/lib/ext/sunec.jar
opt/ums/jre8/lib/ext/sunjce_provider.jar
opt/ums/jre8/lib/ext/sunpkcs11.jar
opt/ums/jre8/lib/ext/zipfs.jar
opt/ums/jre8/lib/flavormap.properties
opt/ums/jre8/lib/hijrah-config-umalqura.properties
opt/ums/jre8/lib/images/
opt/ums/jre8/lib/images/cursors/
opt/ums/jre8/lib/images/cursors/cursors.properties
opt/ums/jre8/lib/images/cursors/invalid32x32.gif
opt/ums/jre8/lib/images/cursors/motif_CopyDrop32x32.gif
opt/ums/jre8/lib/images/cursors/motif_CopyNoDrop32x32.gif
opt/ums/jre8/lib/images/cursors/motif_LinkDrop32x32.gif
opt/ums/jre8/lib/images/cursors/motif_LinkNoDrop32x32.gif
opt/ums/jre8/lib/images/cursors/motif_MoveDrop32x32.gif
opt/ums/jre8/lib/images/cursors/motif_MoveNoDrop32x32.gif
opt/ums/jre8/lib/jce.jar
opt/ums/jre8/lib/jexec
opt/ums/jre8/lib/jfr.jar
opt/ums/jre8/lib/jfr/
opt/ums/jre8/lib/jfr/default.jfc
opt/ums/jre8/lib/jfr/profile.jfc
opt/ums/jre8/lib/jsse.jar
opt/ums/jre8/lib/jvm.hprof.txt
opt/ums/jre8/lib/logging.properties
opt/ums/jre8/lib/management-agent.jar
opt/ums/jre8/lib/management/
opt/ums/jre8/lib/management/jmxremote.access
opt/ums/jre8/lib/management/jmxremote.password.template
opt/ums/jre8/lib/management/management.properties
opt/ums/jre8/lib/management/snmp.acl.template
opt/ums/jre8/lib/meta-index
opt/ums/jre8/lib/net.properties
opt/ums/jre8/lib/psfont.properties.ja
opt/ums/jre8/lib/psfontj2d.properties
opt/ums/jre8/lib/resources.jar
opt/ums/jre8/lib/rt.jar
opt/ums/jre8/lib/security/
opt/ums/jre8/lib/security/blacklisted.certs
opt/ums/jre8/lib/security/cacerts
opt/ums/jre8/lib/security/java.policy
opt/ums/jre8/lib/security/java.security
opt/ums/jre8/lib/security/policy/
opt/ums/jre8/lib/security/policy/limited/
opt/ums/jre8/lib/security/policy/limited/US_export_policy.jar
opt/ums/jre8/lib/security/policy/limited/local_policy.jar
opt/ums/jre8/lib/security/policy/unlimited/
opt/ums/jre8/lib/security/policy/unlimited/US_export_policy.jar
opt/ums/jre8/lib/security/policy/unlimited/local_policy.jar
opt/ums/jre8/lib/sound.properties
opt/ums/jre8/lib/tzdb.dat
opt/ums/jre8/release
opt/ums/linux/
opt/ums/linux/ffmpeg
opt/ums/linux/tsMuxeR
opt/ums/linux/tsMuxeR-new
opt/ums/linux/youtube-dl
opt/ums/logback.headless.xml
opt/ums/logback.xml
opt/ums/plugins/
opt/ums/plugins/README.md
opt/ums/plugins/README.txt
opt/ums/renderers/
opt/ums/renderers/Amazon-FireTVStick-VimuPlayer.conf
opt/ums/renderers/AnyCast.conf
opt/ums/renderers/Apple-TV-4K-VLC.conf
opt/ums/renderers/Apple-TV-VLC.conf
opt/ums/renderers/Apple-iDevice-AirPlayer.conf
opt/ums/renderers/Apple-iDevice-VLC.conf
opt/ums/renderers/Apple-iDevice-VLC32bit.conf
opt/ums/renderers/Apple-iDevice.conf
opt/ums/renderers/BlackBerry-PlayBook-KalemSoftMP.conf
opt/ums/renderers/Bush-FreeviewHighDefinitionDigitalSetTopBox.conf
opt/ums/renderers/CambridgeAudio-AzurBD.conf
opt/ums/renderers/DLink-DSM510.conf
opt/ums/renderers/DefaultRenderer.conf
opt/ums/renderers/DefaultRenderer.conf.orig
opt/ums/renderers/Denon-4311CI.conf
opt/ums/renderers/Denon-X4200W.conf
opt/ums/renderers/DirecTV.conf
opt/ums/renderers/FetchTV.conf
opt/ums/renderers/Free-Freebox.conf
opt/ums/renderers/Freecom-MusicPal.conf
opt/ums/renderers/Google-Android-BubbleUPnP-MXPlayer.conf
opt/ums/renderers/Google-Android-Chromecast.conf
opt/ums/renderers/Google-Android.conf
opt/ums/renderers/Google-ChromecastUltra.conf
opt/ums/renderers/Hama-IR320.conf
opt/ums/renderers/Hisense-K680.conf
opt/ums/renderers/Kodi.conf
opt/ums/renderers/LG-BP.conf
opt/ums/renderers/LG-BP550.conf
opt/ums/renderers/LG-EG910V.conf
opt/ums/renderers/LG-LA6200.conf
opt/ums/renderers/LG-LA644V.conf
opt/ums/renderers/LG-LB.conf
opt/ums/renderers/LG-LM620.conf
opt/ums/renderers/LG-LM660.conf
opt/ums/renderers/LG-LS5700.conf
opt/ums/renderers/LG-OLED.conf
opt/ums/renderers/LG-ST600.conf
opt/ums/renderers/LG-UB820V.conf
opt/ums/renderers/LG-UH770.conf
opt/ums/renderers/LG-WebOS.conf
opt/ums/renderers/Logitech-Squeezebox.conf
opt/ums/renderers/Mediaplayer.conf
opt/ums/renderers/Microsoft-WindowsMediaPlayer.conf
opt/ums/renderers/Microsoft-Xbox360.conf
opt/ums/renderers/Microsoft-XboxOne.conf
opt/ums/renderers/Miracast-M806.conf
opt/ums/renderers/Mirascreen.conf
opt/ums/renderers/Movian.conf
opt/ums/renderers/Naim-Mu-So.conf
opt/ums/renderers/Netgear-NeoTV.conf
opt/ums/renderers/Netgem-N7700.conf
opt/ums/renderers/Nextcp2.conf
opt/ums/renderers/Nokia-N900.conf
opt/ums/renderers/OPPO-BDP.conf
opt/ums/renderers/OPPO-BDP83.conf
opt/ums/renderers/Onkyo-TXNR7xx.conf
opt/ums/renderers/Onkyo-TXNR8xx.conf
opt/ums/renderers/Panasonic-DMPBDT.conf
opt/ums/renderers/Panasonic-DMPBDT220.conf
opt/ums/renderers/Panasonic-DMPBDT360.conf
opt/ums/renderers/Panasonic-DMR.conf
opt/ums/renderers/Panasonic-HZ1500.conf
opt/ums/renderers/Panasonic-SCBTT.conf
opt/ums/renderers/Panasonic-Viera.conf
opt/ums/renderers/Panasonic-VieraAS600E.conf
opt/ums/renderers/Panasonic-VieraAS650.conf
opt/ums/renderers/Panasonic-VieraCX680.conf
opt/ums/renderers/Panasonic-VieraCX700.conf
opt/ums/renderers/Panasonic-VieraDX.conf
opt/ums/renderers/Panasonic-VieraE6.conf
opt/ums/renderers/Panasonic-VieraET60.conf
opt/ums/renderers/Panasonic-VieraGT50.conf
opt/ums/renderers/Panasonic-VieraGX800B.conf
opt/ums/renderers/Panasonic-VieraS60.conf
opt/ums/renderers/Panasonic-VieraST60.conf
opt/ums/renderers/Panasonic-VieraTHPU30Z.conf
opt/ums/renderers/Panasonic-VieraTXL32V10E.conf
opt/ums/renderers/Panasonic-VieraVT60.conf
opt/ums/renderers/Philips-AndroidTV.conf
opt/ums/renderers/Philips-AureaAndNetTV.conf
opt/ums/renderers/Philips-PFL.conf
opt/ums/renderers/Philips-PUS-6500Series.conf
opt/ums/renderers/Philips-PUS.conf
opt/ums/renderers/Philips-Streamium.conf
opt/ums/renderers/Pioneer-BDP.conf
opt/ums/renderers/Pioneer-Kuro.conf
opt/ums/renderers/PopcornHour.conf
opt/ums/renderers/README.txt
opt/ums/renderers/Realtek.conf
opt/ums/renderers/Roku-4640x.conf
opt/ums/renderers/Roku-Roku3-3.conf
opt/ums/renderers/Roku-Roku3-5.conf
opt/ums/renderers/Roku-Roku3-6-7.conf
opt/ums/renderers/Roku-TV8.conf
opt/ums/renderers/Roku-Ultra.conf
opt/ums/renderers/Samsung-8series.conf
opt/ums/renderers/Samsung-9series.conf
opt/ums/renderers/Samsung-BDC6800.conf
opt/ums/renderers/Samsung-BDH6500.conf
opt/ums/renderers/Samsung-C6600.conf
opt/ums/renderers/Samsung-CD.conf
opt/ums/renderers/Samsung-D6400.conf
opt/ums/renderers/Samsung-D7000.conf
opt/ums/renderers/Samsung-EH5300.conf
opt/ums/renderers/Samsung-EH6070.conf
opt/ums/renderers/Samsung-ES6100.conf
opt/ums/renderers/Samsung-ES6575.conf
opt/ums/renderers/Samsung-ES8000.conf
opt/ums/renderers/Samsung-ES8005.conf
opt/ums/renderers/Samsung-F5100.conf
opt/ums/renderers/Samsung-F5505.conf
opt/ums/renderers/Samsung-F5900.conf
opt/ums/renderers/Samsung-GalaxyNoteTab.conf
opt/ums/renderers/Samsung-GalaxyS5.conf
opt/ums/renderers/Samsung-GalaxyS7.conf
opt/ums/renderers/Samsung-H4500.conf
opt/ums/renderers/Samsung-H6203.conf
opt/ums/renderers/Samsung-H6400.conf
opt/ums/renderers/Samsung-HTE3.conf
opt/ums/renderers/Samsung-HTF4.conf
opt/ums/renderers/Samsung-HU7000.conf
opt/ums/renderers/Samsung-HU9000.conf
opt/ums/renderers/Samsung-J55xx.conf
opt/ums/renderers/Samsung-J6200.conf
opt/ums/renderers/Samsung-JU6400.conf
opt/ums/renderers/Samsung-Mobile.conf
opt/ums/renderers/Samsung-NotCD.conf
opt/ums/renderers/Samsung-PL51E490.conf
opt/ums/renderers/Samsung-Q6.conf
opt/ums/renderers/Samsung-Q70.conf
opt/ums/renderers/Samsung-Q9.conf
opt/ums/renderers/Samsung-SMTG7400.conf
opt/ums/renderers/Samsung-Soundbar-MS750.conf
opt/ums/renderers/Samsung-UHD-2019-8K.conf
opt/ums/renderers/Samsung-UHD-2019.conf
opt/ums/renderers/Samsung-UHD.conf
opt/ums/renderers/Samsung-WiseLink.conf
opt/ums/renderers/Sharp-Aquos.conf
opt/ums/renderers/Showtime3.conf
opt/ums/renderers/Showtime4.conf
opt/ums/renderers/Sony-Bluray-BDP-S3700.conf
opt/ums/renderers/Sony-Bluray-UBP-X800M2.conf
opt/ums/renderers/Sony-Bluray.conf
opt/ums/renderers/Sony-Bluray2013.conf
opt/ums/renderers/Sony-Bravia4500.conf
opt/ums/renderers/Sony-Bravia5500.conf
opt/ums/renderers/Sony-BraviaAG.conf
opt/ums/renderers/Sony-BraviaBX305.conf
opt/ums/renderers/Sony-BraviaEX.conf
opt/ums/renderers/Sony-BraviaEX620.conf
opt/ums/renderers/Sony-BraviaEX725.conf
opt/ums/renderers/Sony-BraviaHX.conf
opt/ums/renderers/Sony-BraviaHX75.conf
opt/ums/renderers/Sony-BraviaNX70x.conf
opt/ums/renderers/Sony-BraviaNX800.conf
opt/ums/renderers/Sony-BraviaW.conf
opt/ums/renderers/Sony-BraviaX.conf
opt/ums/renderers/Sony-BraviaXBR-OLED.conf
opt/ums/renderers/Sony-BraviaXBR.conf
opt/ums/renderers/Sony-BraviaXD.conf
opt/ums/renderers/Sony-BraviaXH.conf
opt/ums/renderers/Sony-HomeTheatreSystem.conf
opt/ums/renderers/Sony-PlayStation3.conf
opt/ums/renderers/Sony-PlayStation4.conf
opt/ums/renderers/Sony-PlayStationVita.conf
opt/ums/renderers/Sony-SMPN100.conf
opt/ums/renderers/Sony-STR-DN1080.conf
opt/ums/renderers/Sony-STR5800ES.conf
opt/ums/renderers/Sony-Xperia.conf
opt/ums/renderers/Sony-XperiaZ3.conf
opt/ums/renderers/Technisat-S1Plus.conf
opt/ums/renderers/Telefunken-TV.conf
opt/ums/renderers/Telstra-Tbox.conf
opt/ums/renderers/Thomson-U3.conf
opt/ums/renderers/VLC-for-desktop.conf
opt/ums/renderers/VideoWeb-VideoWebTV.conf
opt/ums/renderers/Vizio-SmartTV.conf
opt/ums/renderers/WesternDigital-WDTVLive.conf
opt/ums/renderers/XBMC.conf
opt/ums/renderers/Yamaha-AV.conf
opt/ums/renderers/Yamaha-RN303.conf
opt/ums/renderers/Yamaha-RN500.conf
opt/ums/renderers/Yamaha-RXA1010.conf
opt/ums/renderers/Yamaha-RXA2050.conf
opt/ums/renderers/Yamaha-RXV3900.conf
opt/ums/renderers/Yamaha-RXV500D.conf
opt/ums/renderers/Yamaha-RXV671.conf
opt/ums/renderers/foobar2000-mobile.conf
opt/ums/ums.jar
opt/ums/web/
opt/ums/web/browse.html
opt/ums/web/browse.html.orig
opt/ums/web/bump/
opt/ums/web/bump/bump.html
opt/ums/web/bump/bump.js
opt/ums/web/bump/bump16-2.png
opt/ums/web/bump/bump32.png
opt/ums/web/bump/skin/
opt/ums/web/bump/skin/add16.png
opt/ums/web/bump/skin/bump16.png
opt/ums/web/bump/skin/clear16.png
opt/ums/web/bump/skin/close16.png
opt/ums/web/bump/skin/fwd16.png
opt/ums/web/bump/skin/mute16.png
opt/ums/web/bump/skin/next16.png
opt/ums/web/bump/skin/pause16.png
opt/ums/web/bump/skin/play16.png
opt/ums/web/bump/skin/prev16.png
opt/ums/web/bump/skin/remove16.png
opt/ums/web/bump/skin/rew16.png
opt/ums/web/bump/skin/skin.js
opt/ums/web/bump/skin/stop16.png
opt/ums/web/bump/skin/vol16.png
opt/ums/web/button-search.png
opt/ums/web/doc.html
opt/ums/web/download.png
opt/ums/web/err32.png
opt/ums/web/favicon.ico
opt/ums/web/flow.html
opt/ums/web/flowplayer.min.js
opt/ums/web/functional.css
opt/ums/web/icon-log.jpg
opt/ums/web/icon-log.png
opt/ums/web/icon-log2.jpg
opt/ums/web/image.html
opt/ums/web/img/
opt/ums/web/img/edit-icon-14.png
opt/ums/web/img/flowplayer.png
opt/ums/web/img/flowplayer@2x.png
opt/ums/web/img/play_white.png
opt/ums/web/img/play_white@x2.png
opt/ums/web/img/play_white_rtl.png
opt/ums/web/img/play_white_rtl@x2.png
opt/ums/web/img/white.png
opt/ums/web/img/white@x2.png
opt/ums/web/img/white_rtl.png
opt/ums/web/img/white_rtl@x2.png
opt/ums/web/info24.png
opt/ums/web/info32.png
opt/ums/web/inputbackground-search.png
opt/ums/web/jquery.min.js
opt/ums/web/jquery.ums.js
opt/ums/web/js.cookie.min.js
opt/ums/web/logo.png
opt/ums/web/logo1.png
opt/ums/web/menu-background-content.png
opt/ums/web/menu-left-content.png
opt/ums/web/ok32.png
opt/ums/web/pad32.png
opt/ums/web/play.html
opt/ums/web/search32.png
opt/ums/web/settings32.png
opt/ums/web/start.html
opt/ums/web/subs.html
opt/ums/web/util/
opt/ums/web/util/bootstrap.min.css
opt/ums/web/util/bootstrap.min.js
opt/ums/web/util/flipper.css
opt/ums/web/util/fontawesome/
opt/ums/web/util/fontawesome/css/
opt/ums/web/util/fontawesome/css/all.min.css
opt/ums/web/util/fontawesome/js/
opt/ums/web/util/fontawesome/js/all.min.js
opt/ums/web/util/fontawesome/webfonts/
opt/ums/web/util/fontawesome/webfonts/fa-brands-400.eot
opt/ums/web/util/fontawesome/webfonts/fa-brands-400.svg
opt/ums/web/util/fontawesome/webfonts/fa-brands-400.ttf
opt/ums/web/util/fontawesome/webfonts/fa-brands-400.woff
opt/ums/web/util/fontawesome/webfonts/fa-brands-400.woff2
opt/ums/web/util/fontawesome/webfonts/fa-regular-400.eot
opt/ums/web/util/fontawesome/webfonts/fa-regular-400.svg
opt/ums/web/util/fontawesome/webfonts/fa-regular-400.ttf
opt/ums/web/util/fontawesome/webfonts/fa-regular-400.woff
opt/ums/web/util/fontawesome/webfonts/fa-regular-400.woff2
opt/ums/web/util/fontawesome/webfonts/fa-solid-900.eot
opt/ums/web/util/fontawesome/webfonts/fa-solid-900.svg
opt/ums/web/util/fontawesome/webfonts/fa-solid-900.ttf
opt/ums/web/util/fontawesome/webfonts/fa-solid-900.woff
opt/ums/web/util/fontawesome/webfonts/fa-solid-900.woff2
opt/ums/web/util/fonts/
opt/ums/web/util/fonts/Roboto-Bold.ttf
opt/ums/web/util/fonts/Roboto-Medium.ttf
opt/ums/web/util/fonts/Roboto-Regular-20140804.ttf
opt/ums/web/util/fonts/Roboto-Regular.ttf
opt/ums/web/util/fonts/codropsicons/
opt/ums/web/util/fonts/codropsicons/codropsicons.eot
opt/ums/web/util/fonts/codropsicons/codropsicons.svg
opt/ums/web/util/fonts/codropsicons/codropsicons.ttf
opt/ums/web/util/fonts/codropsicons/codropsicons.woff
opt/ums/web/util/fonts/codropsicons/license.txt
opt/ums/web/util/fonts/kabuki-material-icons.ttf
opt/ums/web/util/fonts/youtube-icons.ttf
opt/ums/web/util/html5.js
opt/ums/web/util/inflate.js
opt/ums/web/util/jQKeyboard.css
opt/ums/web/util/jQKeyboard.js
opt/ums/web/util/keycode.js
opt/ums/web/util/keycontroller.js
opt/ums/web/util/log.html
opt/ums/web/util/logviewer.html
opt/ums/web/util/modernizr.custom.js
opt/ums/web/util/remoteui.css
opt/ums/web/util/search.js
opt/ums/web/util/searchclassie.js
opt/ums/web/util/sh.css
opt/ums/web/util/sh.js
opt/ums/web/util/shCore.js
opt/ums/web/util/video-js/
opt/ums/web/util/video-js/video-js.min.css
opt/ums/web/util/video-js/video.min.js
opt/ums/web/util/zip.js
opt/ums/web/view-dynamic.png
opt/ums/web/view-grid.png
opt/ums/web/view-info.png
opt/ums/web/view-padcolor.png
opt/ums/web/warn32.png
opt/ums/web/web-narrow.css
opt/ums/web/web-wide.css
opt/ums/web/web.css
usr/
usr/bin/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ums@.service
usr/lib/systemd/system/ums@.timer
usr/share/
usr/share/applications/
usr/share/applications/ums.desktop
usr/share/pixmaps/
usr/share/pixmaps/ums.png
