# files/gerrit-3.4.1-1/files
%FILES%
etc/
etc/default/
etc/default/gerritcodereview
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/gerrit.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/gerrit.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/gerrit.conf
usr/share/
usr/share/java/
usr/share/java/gerrit/
usr/share/java/gerrit/gerrit.war
