# files/openssh-hpn-9.0p1.hpn17v2-1/files
%FILES%
etc/
etc/hpnssh/
etc/hpnssh/moduli
etc/hpnssh/ssh_config
etc/hpnssh/sshd_config
etc/pam.d/
etc/pam.d/hpnsshd
usr/
usr/bin/
usr/bin/hpnfindssl.sh
usr/bin/hpnscp
usr/bin/hpnsftp
usr/bin/hpnssh
usr/bin/hpnssh-add
usr/bin/hpnssh-agent
usr/bin/hpnssh-copy-id
usr/bin/hpnssh-keygen
usr/bin/hpnssh-keyscan
usr/bin/hpnsshd
usr/lib/
usr/lib/hpnssh/
usr/lib/hpnssh/hpnsftp-server
usr/lib/hpnssh/hpnssh-keysign
usr/lib/hpnssh/hpnssh-pkcs11-helper
usr/lib/hpnssh/hpnssh-sk-helper
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hpnsshd.service
usr/lib/systemd/system/hpnsshdgenkeys.service
usr/share/
usr/share/licenses/
usr/share/licenses/openssh-hpn/
usr/share/licenses/openssh-hpn/LICENCE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/hpnscp.1.gz
usr/share/man/man1/hpnsftp.1.gz
usr/share/man/man1/hpnssh-add.1.gz
usr/share/man/man1/hpnssh-agent.1.gz
usr/share/man/man1/hpnssh-copy-id.1.gz
usr/share/man/man1/hpnssh-keygen.1.gz
usr/share/man/man1/hpnssh-keyscan.1.gz
usr/share/man/man1/hpnssh.1.gz
usr/share/man/man1/slogin.1.gz
usr/share/man/man5/
usr/share/man/man5/hpnmoduli.5.gz
usr/share/man/man5/hpnssh_config.5.gz
usr/share/man/man5/hpnsshd_config.5.gz
usr/share/man/man8/
usr/share/man/man8/hpnsftp-server.8.gz
usr/share/man/man8/hpnssh-keysign.8.gz
usr/share/man/man8/hpnssh-pkcs11-helper.8.gz
usr/share/man/man8/hpnssh-sk-helper.8.gz
usr/share/man/man8/hpnsshd.8.gz
var/
var/empty/
var/lib/
var/lib/tmpfiles.d/
var/lib/tmpfiles.d/hpnsshd.conf
