# files/sysbox-ce-0.4.1-1/files
%FILES%
usr/
usr/bin/
usr/bin/sysbox
usr/bin/sysbox-fs
usr/bin/sysbox-mgr
usr/bin/sysbox-runc
usr/lib/
usr/lib/modules-load.d/
usr/lib/modules-load.d/50-sysbox-mod.conf
usr/lib/sysctl.d/
usr/lib/sysctl.d/99-sysbox-sysctl.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sysbox-fs.service
usr/lib/systemd/system/sysbox-mgr.service
usr/lib/systemd/system/sysbox.service
usr/share/
usr/share/doc/
usr/share/doc/sysbox-ce/
usr/share/doc/sysbox-ce/changelog.Debian.gz
usr/share/doc/sysbox-ce/copyright
