# files/x2goserver-4.1.0.3-10/files
%FILES%
etc/
etc/sudoers.d/
etc/sudoers.d/x2goserver
etc/x2go/
etc/x2go/desktopsharing/
etc/x2go/desktopsharing/settings
etc/x2go/keystrokes.cfg
etc/x2go/x2go_logout
etc/x2go/x2go_logout.d/
etc/x2go/x2go_logout.d/010_userscripts.sh
etc/x2go/x2goagent.keyboard
etc/x2go/x2goagent.options
etc/x2go/x2goserver.conf
etc/x2go/x2gosql/
etc/x2go/x2gosql/passwords/
etc/x2go/x2gosql/sql
usr/
usr/bin/
usr/bin/x2goagent
usr/bin/x2gobasepath
usr/bin/x2gocleansessions
usr/bin/x2gocmdexitmessage
usr/bin/x2godbadmin
usr/bin/x2gofeature
usr/bin/x2gofeaturelist
usr/bin/x2gofm
usr/bin/x2gogetapps
usr/bin/x2gogetservers
usr/bin/x2golistdesktops
usr/bin/x2golistmounts
usr/bin/x2golistsessions
usr/bin/x2golistsessions_root
usr/bin/x2golistshadowsessions
usr/bin/x2gomountdirs
usr/bin/x2gopath
usr/bin/x2goprint
usr/bin/x2goresume-desktopsharing
usr/bin/x2goresume-session
usr/bin/x2goruncommand
usr/bin/x2goserver-run-extensions
usr/bin/x2gosessionlimit
usr/bin/x2gosetkeyboard
usr/bin/x2goshowblocks
usr/bin/x2gostartagent
usr/bin/x2gosuspend-desktopsharing
usr/bin/x2gosuspend-session
usr/bin/x2goterminate-desktopsharing
usr/bin/x2goterminate-session
usr/bin/x2goumount-session
usr/bin/x2goversion
usr/lib/
usr/lib/modules-load.d/
usr/lib/modules-load.d/x2goserver.conf
usr/lib/nx/
usr/lib/nx/bin/
usr/lib/nx/bin/x2goagent
usr/lib/perl5/
usr/lib/perl5/5.34/
usr/lib/perl5/5.34/vendor_perl/
usr/lib/perl5/5.34/vendor_perl/auto/
usr/lib/perl5/5.34/vendor_perl/auto/x2goserver/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/x2goserver.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/x2goserver.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/x2goserver.conf
usr/lib/x2go/
usr/lib/x2go/extensions/
usr/lib/x2go/extensions/fail-resume.d/
usr/lib/x2go/extensions/fail-resume.d/.placeholder
usr/lib/x2go/extensions/fail-runcommand.d/
usr/lib/x2go/extensions/fail-runcommand.d/.placeholder
usr/lib/x2go/extensions/fail-start.d/
usr/lib/x2go/extensions/fail-start.d/.placeholder
usr/lib/x2go/extensions/fail-suspend.d/
usr/lib/x2go/extensions/fail-suspend.d/.placeholder
usr/lib/x2go/extensions/fail-terminate.d/
usr/lib/x2go/extensions/fail-terminate.d/.placeholder
usr/lib/x2go/extensions/post-resume.d/
usr/lib/x2go/extensions/post-resume.d/.placeholder
usr/lib/x2go/extensions/post-runcommand.d/
usr/lib/x2go/extensions/post-runcommand.d/.placeholder
usr/lib/x2go/extensions/post-start.d/
usr/lib/x2go/extensions/post-start.d/.placeholder
usr/lib/x2go/extensions/post-suspend.d/
usr/lib/x2go/extensions/post-suspend.d/.placeholder
usr/lib/x2go/extensions/post-terminate.d/
usr/lib/x2go/extensions/post-terminate.d/.placeholder
usr/lib/x2go/extensions/pre-resume.d/
usr/lib/x2go/extensions/pre-resume.d/.placeholder
usr/lib/x2go/extensions/pre-runcommand.d/
usr/lib/x2go/extensions/pre-runcommand.d/.placeholder
usr/lib/x2go/extensions/pre-start.d/
usr/lib/x2go/extensions/pre-start.d/.placeholder
usr/lib/x2go/extensions/pre-suspend.d/
usr/lib/x2go/extensions/pre-suspend.d/.placeholder
usr/lib/x2go/extensions/pre-terminate.d/
usr/lib/x2go/extensions/pre-terminate.d/.placeholder
usr/lib/x2go/libx2go-server-db-sqlite3-wrapper
usr/lib/x2go/libx2go-server-db-sqlite3-wrapper.pl
usr/lib/x2go/x2gochangestatus
usr/lib/x2go/x2gocheckport
usr/lib/x2go/x2gocreatesession
usr/lib/x2go/x2gocreateshadowsession
usr/lib/x2go/x2gogetagent
usr/lib/x2go/x2gogetagentstate
usr/lib/x2go/x2gogetdisplays
usr/lib/x2go/x2gogetfreeport
usr/lib/x2go/x2gogetports
usr/lib/x2go/x2gogetrandomport
usr/lib/x2go/x2gogetstatus
usr/lib/x2go/x2goinsertport
usr/lib/x2go/x2goinsertsession
usr/lib/x2go/x2goinsertshadowsession
usr/lib/x2go/x2goistrue
usr/lib/x2go/x2golistsessions_sql
usr/lib/x2go/x2gologlevel
usr/lib/x2go/x2goqueryconfig
usr/lib/x2go/x2goresume
usr/lib/x2go/x2gormforward
usr/lib/x2go/x2gormport
usr/lib/x2go/x2gosuspend-agent
usr/lib/x2go/x2gosyslog
usr/share/
usr/share/applications/
usr/share/applications/x2gofm.desktop
usr/share/doc/
usr/share/doc/x2goserver/
usr/share/doc/x2goserver/ChangeLog
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/x2goagent.1.gz
usr/share/man/man3/
usr/share/man/man3/X2Go::Config.3pm.gz
usr/share/man/man3/X2Go::Log.3pm.gz
usr/share/man/man3/X2Go::Server.3pm.gz
usr/share/man/man3/X2Go::Server::Agent.3pm.gz
usr/share/man/man3/X2Go::Server::Agent::NX.3pm.gz
usr/share/man/man3/X2Go::Server::DB.3pm.gz
usr/share/man/man3/X2Go::Server::DB::PostgreSQL.3pm.gz
usr/share/man/man3/X2Go::Server::DB::SQLite3.3pm.gz
usr/share/man/man3/X2Go::SupeReNicer.3pm.gz
usr/share/man/man3/X2Go::Utils.3pm.gz
usr/share/man/man5/
usr/share/man/man5/x2goserver.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/x2gobasepath.8.gz
usr/share/man/man8/x2gocmdexitmessage.8.gz
usr/share/man/man8/x2godbadmin.8.gz
usr/share/man/man8/x2gofeature.8.gz
usr/share/man/man8/x2gofeaturelist.8.gz
usr/share/man/man8/x2gofm.8.gz
usr/share/man/man8/x2gogetapps.8.gz
usr/share/man/man8/x2gogetservers.8.gz
usr/share/man/man8/x2golistdesktops.8.gz
usr/share/man/man8/x2golistmounts.8.gz
usr/share/man/man8/x2golistsessions.8.gz
usr/share/man/man8/x2golistsessions_root.8.gz
usr/share/man/man8/x2golistshadowsessions.8.gz
usr/share/man/man8/x2gomountdirs.8.gz
usr/share/man/man8/x2gopath.8.gz
usr/share/man/man8/x2goprint.8.gz
usr/share/man/man8/x2goresume-desktopsharing.8.gz
usr/share/man/man8/x2goresume-session.8.gz
usr/share/man/man8/x2goruncommand.8.gz
usr/share/man/man8/x2goserver-run-extensions.8.gz
usr/share/man/man8/x2gosessionlimit.8.gz
usr/share/man/man8/x2gosetkeyboard.8.gz
usr/share/man/man8/x2goshowblocks.8.gz
usr/share/man/man8/x2gostartagent.8.gz
usr/share/man/man8/x2gosuspend-desktopsharing.8.gz
usr/share/man/man8/x2gosuspend-session.8.gz
usr/share/man/man8/x2goterminate-desktopsharing.8.gz
usr/share/man/man8/x2goterminate-session.8.gz
usr/share/man/man8/x2goumount-session.8.gz
usr/share/man/man8/x2goversion.8.gz
usr/share/mime/
usr/share/mime/packages/
usr/share/mime/packages/sshfs-x2go.xml
usr/share/perl5/
usr/share/perl5/vendor_perl/
usr/share/perl5/vendor_perl/X2Go/
usr/share/perl5/vendor_perl/X2Go/Config.pm
usr/share/perl5/vendor_perl/X2Go/Log.pm
usr/share/perl5/vendor_perl/X2Go/Server.pm
usr/share/perl5/vendor_perl/X2Go/Server/
usr/share/perl5/vendor_perl/X2Go/Server/Agent.pm
usr/share/perl5/vendor_perl/X2Go/Server/Agent/
usr/share/perl5/vendor_perl/X2Go/Server/Agent/NX.pm
usr/share/perl5/vendor_perl/X2Go/Server/DB.pm
usr/share/perl5/vendor_perl/X2Go/Server/DB/
usr/share/perl5/vendor_perl/X2Go/Server/DB/PostgreSQL.pm
usr/share/perl5/vendor_perl/X2Go/Server/DB/SQLite3.pm
usr/share/perl5/vendor_perl/X2Go/SupeReNicer.pm
usr/share/perl5/vendor_perl/X2Go/Utils.pm
usr/share/pixmaps/
usr/share/pixmaps/x2goagent.xpm
usr/share/x2go/
usr/share/x2go/versions/
usr/share/x2go/versions/VERSION.x2goserver
usr/share/x2go/versions/VERSION.x2goserver-common
usr/share/x2go/versions/VERSION.x2goserver-desktopsharing
usr/share/x2go/versions/VERSION.x2goserver-extensions
usr/share/x2go/versions/VERSION.x2goserver-fmbindings
usr/share/x2go/versions/VERSION.x2goserver-printing
usr/share/x2go/versions/VERSION.x2goserver-x2goagent
usr/share/x2go/x2gofeature.d/
usr/share/x2go/x2gofeature.d/x2goserver-desktopsharing.features
usr/share/x2go/x2gofeature.d/x2goserver-extensions.features
usr/share/x2go/x2gofeature.d/x2goserver-fmbindings.features
usr/share/x2go/x2gofeature.d/x2goserver-printing.features
usr/share/x2go/x2gofeature.d/x2goserver-x2goagent.features
usr/share/x2go/x2gofeature.d/x2goserver.features
