# files/switcheroo-control-2.6-1/files
%FILES%
usr/
usr/bin/
usr/bin/switcherooctl
usr/lib/
usr/lib/switcheroo-control
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/switcheroo-control.service
usr/lib/udev/
usr/lib/udev/hwdb.d/
usr/lib/udev/hwdb.d/30-pci-intel-gpu.hwdb
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/net.hadess.SwitcherooControl.conf
