# files/freenet-0.7.5.1494-1/files
%FILES%
opt/
opt/freenet/
opt/freenet/conf/
opt/freenet/conf/freenet.ini.dist
opt/freenet/downloads/
opt/freenet/lib/
opt/freenet/lib/bcprov-jdk18on-1.71.jar
opt/freenet/lib/freenet-ext-29.jar
opt/freenet/lib/freenet.jar
opt/freenet/lib/jna-4.5.2.jar
opt/freenet/lib/jna-platform-4.5.2.jar
opt/freenet/noderef/
opt/freenet/noderef/seednodes.fref
opt/freenet/persistent-temp/
opt/freenet/plugins/
opt/freenet/plugins/WebOfTrust.jar
opt/freenet/plugins/data/
opt/freenet/tmp/
opt/freenet/user/
opt/freenet/user/certs/
opt/freenet/user/data/
opt/freenet/wrapper.config
run/
run/freenet/
usr/
usr/bin/
usr/bin/freenet
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/freenet.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/freenet.conf
usr/share/
usr/share/licenses/
usr/share/licenses/freenet/
usr/share/licenses/freenet/LICENSE
