# files/prelockd-git-0.9.r4.g2c408bd-1/files
%FILES%
etc/
etc/prelockd.conf
usr/
usr/bin/
usr/bin/prelockd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prelockd.service
usr/share/
usr/share/doc/
usr/share/doc/prelockd/
usr/share/doc/prelockd/MANPAGE.md
usr/share/doc/prelockd/README.md
usr/share/licenses/
usr/share/licenses/prelockd-git/
usr/share/licenses/prelockd-git/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/prelockd.8.gz
usr/share/prelockd/
usr/share/prelockd/prelockd.conf
var/
var/lib/
var/lib/prelockd/
