# files/connman-git-1.41.r34.g1e1d2af0-1/files
%FILES%
usr/
usr/bin/
usr/bin/connman-vpnd
usr/bin/connmanctl
usr/bin/connmand
usr/bin/connmand-wait-online
usr/include/
usr/include/connman/
usr/include/connman/acd.h
usr/include/connman/agent.h
usr/include/connman/device.h
usr/include/connman/inet.h
usr/include/connman/inotify.h
usr/include/connman/ipaddress.h
usr/include/connman/ipconfig.h
usr/include/connman/log.h
usr/include/connman/machine.h
usr/include/connman/network.h
usr/include/connman/notifier.h
usr/include/connman/peer.h
usr/include/connman/plugin.h
usr/include/connman/provision.h
usr/include/connman/resolver.h
usr/include/connman/service.h
usr/include/connman/session.h
usr/include/connman/storage.h
usr/include/connman/tethering.h
usr/include/connman/version.h
usr/lib/
usr/lib/connman/
usr/lib/connman/plugins-vpn/
usr/lib/connman/plugins-vpn/openconnect.so
usr/lib/connman/plugins-vpn/openvpn.so
usr/lib/connman/plugins-vpn/pptp.so
usr/lib/connman/plugins-vpn/vpnc.so
usr/lib/connman/plugins-vpn/wireguard.so
usr/lib/connman/scripts/
usr/lib/connman/scripts/libppp-plugin.so
usr/lib/connman/scripts/openvpn-script
usr/lib/connman/scripts/vpn-script
usr/lib/pkgconfig/
usr/lib/pkgconfig/connman.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/connman-vpn.service
usr/lib/systemd/system/connman-wait-online.service
usr/lib/systemd/system/connman.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/connman_resolvconf.conf
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/net.connman.vpn.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/connman-nmcompat.conf
usr/share/dbus-1/system.d/connman-vpn-dbus.conf
usr/share/dbus-1/system.d/connman.conf
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/connmanctl.1.gz
usr/share/man/man5/
usr/share/man/man5/connman-service.config.5.gz
usr/share/man/man5/connman-vpn-provider.config.5.gz
usr/share/man/man5/connman-vpn.conf.5.gz
usr/share/man/man5/connman.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/connman-vpn.8.gz
usr/share/man/man8/connman.8.gz
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/net.connman.policy
usr/share/polkit-1/actions/net.connman.vpn.policy
