# files/sonarr-3.0.6.1342-1/files
%FILES%
usr/
usr/lib/
usr/lib/sonarr/
usr/lib/sonarr/bin/
usr/lib/sonarr/bin/0Harmony.dll
usr/lib/sonarr/bin/BouncyCastle.Crypto.dll
usr/lib/sonarr/bin/CookComputing.XmlRpcV2.dll
usr/lib/sonarr/bin/FluentMigrator.Runner.dll
usr/lib/sonarr/bin/FluentMigrator.dll
usr/lib/sonarr/bin/FluentValidation.dll
usr/lib/sonarr/bin/ICSharpCode.SharpZipLib.dll
usr/lib/sonarr/bin/Ical.Net.Collections.dll
usr/lib/sonarr/bin/Ical.Net.dll
usr/lib/sonarr/bin/LogentriesCore.dll
usr/lib/sonarr/bin/LogentriesCore.dll.mdb
usr/lib/sonarr/bin/LogentriesNLog.dll
usr/lib/sonarr/bin/LogentriesNLog.dll.mdb
usr/lib/sonarr/bin/MailKit.dll
usr/lib/sonarr/bin/Marr.Data.dll
usr/lib/sonarr/bin/Marr.Data.dll.mdb
usr/lib/sonarr/bin/Microsoft.AspNet.SignalR.Core.dll
usr/lib/sonarr/bin/Microsoft.AspNet.SignalR.SystemWeb.dll
usr/lib/sonarr/bin/Microsoft.Owin.Diagnostics.dll
usr/lib/sonarr/bin/Microsoft.Owin.Host.HttpListener.dll
usr/lib/sonarr/bin/Microsoft.Owin.Host.SystemWeb.dll
usr/lib/sonarr/bin/Microsoft.Owin.Hosting.dll
usr/lib/sonarr/bin/Microsoft.Owin.Security.dll
usr/lib/sonarr/bin/Microsoft.Owin.dll
usr/lib/sonarr/bin/Microsoft.Win32.Primitives.dll
usr/lib/sonarr/bin/MimeKit.dll
usr/lib/sonarr/bin/MonoTorrent.dll
usr/lib/sonarr/bin/MonoTorrent.dll.mdb
usr/lib/sonarr/bin/NLog.dll
usr/lib/sonarr/bin/Nancy.Authentication.Basic.dll
usr/lib/sonarr/bin/Nancy.Authentication.Forms.dll
usr/lib/sonarr/bin/Nancy.Owin.dll
usr/lib/sonarr/bin/Nancy.dll
usr/lib/sonarr/bin/Newtonsoft.Json.dll
usr/lib/sonarr/bin/NodaTime.dll
usr/lib/sonarr/bin/OAuth.dll
usr/lib/sonarr/bin/Org.Mentalis.dll
usr/lib/sonarr/bin/Owin.dll
usr/lib/sonarr/bin/RestSharp.dll
usr/lib/sonarr/bin/Sentry.PlatformAbstractions.dll
usr/lib/sonarr/bin/Sentry.Protocol.dll
usr/lib/sonarr/bin/Sentry.dll
usr/lib/sonarr/bin/SixLabors.Core.dll
usr/lib/sonarr/bin/SixLabors.Fonts.dll
usr/lib/sonarr/bin/SixLabors.ImageSharp.Drawing.dll
usr/lib/sonarr/bin/SixLabors.ImageSharp.dll
usr/lib/sonarr/bin/SixLabors.Shapes.dll
usr/lib/sonarr/bin/SocksWebProxy.dll
usr/lib/sonarr/bin/Sonarr.Api.V3.dll
usr/lib/sonarr/bin/Sonarr.Api.V3.dll.mdb
usr/lib/sonarr/bin/Sonarr.Api.dll
usr/lib/sonarr/bin/Sonarr.Api.dll.mdb
usr/lib/sonarr/bin/Sonarr.Common.dll
usr/lib/sonarr/bin/Sonarr.Common.dll.mdb
usr/lib/sonarr/bin/Sonarr.Core.dll
usr/lib/sonarr/bin/Sonarr.Core.dll.config
usr/lib/sonarr/bin/Sonarr.Core.dll.mdb
usr/lib/sonarr/bin/Sonarr.Host.dll
usr/lib/sonarr/bin/Sonarr.Host.dll.mdb
usr/lib/sonarr/bin/Sonarr.Http.dll
usr/lib/sonarr/bin/Sonarr.Http.dll.mdb
usr/lib/sonarr/bin/Sonarr.Mono.dll
usr/lib/sonarr/bin/Sonarr.Mono.dll.mdb
usr/lib/sonarr/bin/Sonarr.RuntimePatches.dll
usr/lib/sonarr/bin/Sonarr.RuntimePatches.dll.mdb
usr/lib/sonarr/bin/Sonarr.SignalR.dll
usr/lib/sonarr/bin/Sonarr.SignalR.dll.mdb
usr/lib/sonarr/bin/Sonarr.exe
usr/lib/sonarr/bin/Sonarr.exe.config
usr/lib/sonarr/bin/Sonarr.exe.mdb
usr/lib/sonarr/bin/System.AppContext.dll
usr/lib/sonarr/bin/System.Buffers.dll
usr/lib/sonarr/bin/System.Collections.Concurrent.dll
usr/lib/sonarr/bin/System.Collections.Immutable.dll
usr/lib/sonarr/bin/System.Collections.NonGeneric.dll
usr/lib/sonarr/bin/System.Collections.Specialized.dll
usr/lib/sonarr/bin/System.Collections.dll
usr/lib/sonarr/bin/System.ComponentModel.Annotations.dll
usr/lib/sonarr/bin/System.ComponentModel.EventBasedAsync.dll
usr/lib/sonarr/bin/System.ComponentModel.Primitives.dll
usr/lib/sonarr/bin/System.ComponentModel.TypeConverter.dll
usr/lib/sonarr/bin/System.ComponentModel.dll
usr/lib/sonarr/bin/System.Console.dll
usr/lib/sonarr/bin/System.Data.Common.dll
usr/lib/sonarr/bin/System.Data.SQLite.dll
usr/lib/sonarr/bin/System.Diagnostics.Contracts.dll
usr/lib/sonarr/bin/System.Diagnostics.Debug.dll
usr/lib/sonarr/bin/System.Diagnostics.FileVersionInfo.dll
usr/lib/sonarr/bin/System.Diagnostics.Process.dll
usr/lib/sonarr/bin/System.Diagnostics.StackTrace.dll
usr/lib/sonarr/bin/System.Diagnostics.TextWriterTraceListener.dll
usr/lib/sonarr/bin/System.Diagnostics.Tools.dll
usr/lib/sonarr/bin/System.Diagnostics.TraceSource.dll
usr/lib/sonarr/bin/System.Diagnostics.Tracing.dll
usr/lib/sonarr/bin/System.Drawing.Primitives.dll
usr/lib/sonarr/bin/System.Dynamic.Runtime.dll
usr/lib/sonarr/bin/System.Globalization.Calendars.dll
usr/lib/sonarr/bin/System.Globalization.Extensions.dll
usr/lib/sonarr/bin/System.Globalization.dll
usr/lib/sonarr/bin/System.IO.Compression.ZipFile.dll
usr/lib/sonarr/bin/System.IO.FileSystem.DriveInfo.dll
usr/lib/sonarr/bin/System.IO.FileSystem.Primitives.dll
usr/lib/sonarr/bin/System.IO.FileSystem.Watcher.dll
usr/lib/sonarr/bin/System.IO.FileSystem.dll
usr/lib/sonarr/bin/System.IO.IsolatedStorage.dll
usr/lib/sonarr/bin/System.IO.MemoryMappedFiles.dll
usr/lib/sonarr/bin/System.IO.Pipes.dll
usr/lib/sonarr/bin/System.IO.UnmanagedMemoryStream.dll
usr/lib/sonarr/bin/System.IO.dll
usr/lib/sonarr/bin/System.Linq.Expressions.dll
usr/lib/sonarr/bin/System.Linq.Parallel.dll
usr/lib/sonarr/bin/System.Linq.Queryable.dll
usr/lib/sonarr/bin/System.Linq.dll
usr/lib/sonarr/bin/System.Memory.dll
usr/lib/sonarr/bin/System.Net.NameResolution.dll
usr/lib/sonarr/bin/System.Net.NetworkInformation.dll
usr/lib/sonarr/bin/System.Net.Ping.dll
usr/lib/sonarr/bin/System.Net.Primitives.dll
usr/lib/sonarr/bin/System.Net.Requests.dll
usr/lib/sonarr/bin/System.Net.Security.dll
usr/lib/sonarr/bin/System.Net.Sockets.dll
usr/lib/sonarr/bin/System.Net.WebHeaderCollection.dll
usr/lib/sonarr/bin/System.Net.WebSockets.Client.dll
usr/lib/sonarr/bin/System.Net.WebSockets.dll
usr/lib/sonarr/bin/System.Numerics.Vectors.dll
usr/lib/sonarr/bin/System.ObjectModel.dll
usr/lib/sonarr/bin/System.Reflection.Extensions.dll
usr/lib/sonarr/bin/System.Reflection.Primitives.dll
usr/lib/sonarr/bin/System.Reflection.dll
usr/lib/sonarr/bin/System.Resources.Reader.dll
usr/lib/sonarr/bin/System.Resources.ResourceManager.dll
usr/lib/sonarr/bin/System.Resources.Writer.dll
usr/lib/sonarr/bin/System.Runtime.CompilerServices.Unsafe.dll
usr/lib/sonarr/bin/System.Runtime.CompilerServices.VisualC.dll
usr/lib/sonarr/bin/System.Runtime.Extensions.dll
usr/lib/sonarr/bin/System.Runtime.Handles.dll
usr/lib/sonarr/bin/System.Runtime.InteropServices.RuntimeInformation.dll
usr/lib/sonarr/bin/System.Runtime.InteropServices.dll
usr/lib/sonarr/bin/System.Runtime.Numerics.dll
usr/lib/sonarr/bin/System.Runtime.Serialization.Formatters.dll
usr/lib/sonarr/bin/System.Runtime.Serialization.Json.dll
usr/lib/sonarr/bin/System.Runtime.Serialization.Primitives.dll
usr/lib/sonarr/bin/System.Runtime.Serialization.Xml.dll
usr/lib/sonarr/bin/System.Runtime.dll
usr/lib/sonarr/bin/System.Security.Claims.dll
usr/lib/sonarr/bin/System.Security.Cryptography.Algorithms.dll
usr/lib/sonarr/bin/System.Security.Cryptography.Csp.dll
usr/lib/sonarr/bin/System.Security.Cryptography.Encoding.dll
usr/lib/sonarr/bin/System.Security.Cryptography.Primitives.dll
usr/lib/sonarr/bin/System.Security.Cryptography.X509Certificates.dll
usr/lib/sonarr/bin/System.Security.Principal.dll
usr/lib/sonarr/bin/System.Security.SecureString.dll
usr/lib/sonarr/bin/System.Text.Encoding.Extensions.dll
usr/lib/sonarr/bin/System.Text.Encoding.dll
usr/lib/sonarr/bin/System.Text.RegularExpressions.dll
usr/lib/sonarr/bin/System.Threading.Overlapped.dll
usr/lib/sonarr/bin/System.Threading.Tasks.Parallel.dll
usr/lib/sonarr/bin/System.Threading.Tasks.dll
usr/lib/sonarr/bin/System.Threading.Thread.dll
usr/lib/sonarr/bin/System.Threading.ThreadPool.dll
usr/lib/sonarr/bin/System.Threading.Timer.dll
usr/lib/sonarr/bin/System.Threading.dll
usr/lib/sonarr/bin/System.ValueTuple.dll
usr/lib/sonarr/bin/System.Xml.ReaderWriter.dll
usr/lib/sonarr/bin/System.Xml.XDocument.dll
usr/lib/sonarr/bin/System.Xml.XPath.XDocument.dll
usr/lib/sonarr/bin/System.Xml.XPath.dll
usr/lib/sonarr/bin/System.Xml.XmlDocument.dll
usr/lib/sonarr/bin/System.Xml.XmlSerializer.dll
usr/lib/sonarr/bin/UI/
usr/lib/sonarr/bin/UI/Content/
usr/lib/sonarr/bin/UI/Content/Fonts/
usr/lib/sonarr/bin/UI/Content/Fonts/Roboto-Light.ttf
usr/lib/sonarr/bin/UI/Content/Fonts/Roboto-Light.woff
usr/lib/sonarr/bin/UI/Content/Fonts/Roboto-Light.woff2
usr/lib/sonarr/bin/UI/Content/Fonts/Roboto-Regular.ttf
usr/lib/sonarr/bin/UI/Content/Fonts/Roboto-Regular.woff
usr/lib/sonarr/bin/UI/Content/Fonts/Roboto-Regular.woff2
usr/lib/sonarr/bin/UI/Content/Fonts/UbuntuMono-Regular.eot
usr/lib/sonarr/bin/UI/Content/Fonts/UbuntuMono-Regular.ttf
usr/lib/sonarr/bin/UI/Content/Fonts/UbuntuMono-Regular.woff
usr/lib/sonarr/bin/UI/Content/Fonts/fonts.css
usr/lib/sonarr/bin/UI/Content/Fonts/text-security-disc.ttf
usr/lib/sonarr/bin/UI/Content/Fonts/text-security-disc.woff
usr/lib/sonarr/bin/UI/Content/Images/
usr/lib/sonarr/bin/UI/Content/Images/404.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/
usr/lib/sonarr/bin/UI/Content/Images/Icons/android-chrome-192x192.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/android-chrome-512x512.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/apple-touch-icon.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/browserconfig.xml
usr/lib/sonarr/bin/UI/Content/Images/Icons/favicon-16x16.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/favicon-32x32.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/favicon-debug-16x16.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/favicon-debug-32x32.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/favicon-debug.ico
usr/lib/sonarr/bin/UI/Content/Images/Icons/favicon.ico
usr/lib/sonarr/bin/UI/Content/Images/Icons/manifest.json
usr/lib/sonarr/bin/UI/Content/Images/Icons/mstile-144x144.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/mstile-150x150.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/mstile-310x150.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/mstile-310x310.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/mstile-70x70.png
usr/lib/sonarr/bin/UI/Content/Images/Icons/safari-pinned-tab.svg
usr/lib/sonarr/bin/UI/Content/Images/error.png
usr/lib/sonarr/bin/UI/Content/Images/logo.svg
usr/lib/sonarr/bin/UI/Content/Images/poster-dark.png
usr/lib/sonarr/bin/UI/Content/robots.txt
usr/lib/sonarr/bin/UI/Content/styles.css
usr/lib/sonarr/bin/UI/Content/styles.css.map
usr/lib/sonarr/bin/UI/index.html
usr/lib/sonarr/bin/UI/index.js
usr/lib/sonarr/bin/UI/index.js.map
usr/lib/sonarr/bin/UI/login.html
usr/lib/sonarr/bin/UI/oauth.html
usr/lib/sonarr/bin/UI/vendors.js
usr/lib/sonarr/bin/UI/vendors.js.LICENSE.txt
usr/lib/sonarr/bin/UI/vendors.js.map
usr/lib/sonarr/bin/antlr.runtime.dll
usr/lib/sonarr/bin/netstandard.dll
usr/lib/sonarr/bin/release_info
usr/lib/sonarr/package_info
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sonarr.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/sonarr.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/sonarr.conf
