# files/jackman-git-r93.e49a782-1/files
%FILES%
etc/
etc/udev/
etc/udev/rules.d/
etc/udev/rules.d/50-jackman.rules
etc/xdg/
etc/xdg/autostart/
usr/
usr/bin/
usr/bin/jackman
usr/bin/jackman_udev_plug
usr/bin/jackman_udev_unplug
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/jackman_plug@.service
usr/lib/systemd/system/jackman_unplug@.service
usr/share/
usr/share/locale/
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/jackman.mo
