# files/pacvis-git-0.2.7.r12.g34f7494-1/files
%FILES%
usr/
usr/bin/
usr/bin/pacvis
usr/lib/
usr/lib/python3.10/
usr/lib/python3.10/site-packages/
usr/lib/python3.10/site-packages/PacVis-0.2.7-py3.10.egg-info/
usr/lib/python3.10/site-packages/PacVis-0.2.7-py3.10.egg-info/PKG-INFO
usr/lib/python3.10/site-packages/PacVis-0.2.7-py3.10.egg-info/SOURCES.txt
usr/lib/python3.10/site-packages/PacVis-0.2.7-py3.10.egg-info/dependency_links.txt
usr/lib/python3.10/site-packages/PacVis-0.2.7-py3.10.egg-info/entry_points.txt
usr/lib/python3.10/site-packages/PacVis-0.2.7-py3.10.egg-info/requires.txt
usr/lib/python3.10/site-packages/PacVis-0.2.7-py3.10.egg-info/top_level.txt
usr/lib/python3.10/site-packages/pacvis/
usr/lib/python3.10/site-packages/pacvis/__init__.py
usr/lib/python3.10/site-packages/pacvis/__pycache__/
usr/lib/python3.10/site-packages/pacvis/__pycache__/__init__.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/pacvis/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/pacvis/__pycache__/console.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/pacvis/__pycache__/console.cpython-310.pyc
usr/lib/python3.10/site-packages/pacvis/__pycache__/infos.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/pacvis/__pycache__/infos.cpython-310.pyc
usr/lib/python3.10/site-packages/pacvis/__pycache__/pacvis.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/pacvis/__pycache__/pacvis.cpython-310.pyc
usr/lib/python3.10/site-packages/pacvis/console.py
usr/lib/python3.10/site-packages/pacvis/infos.py
usr/lib/python3.10/site-packages/pacvis/pacvis.py
usr/lib/python3.10/site-packages/pacvis/static/
usr/lib/python3.10/site-packages/pacvis/static/MaterialIcons-Regular.eot
usr/lib/python3.10/site-packages/pacvis/static/MaterialIcons-Regular.ijmap
usr/lib/python3.10/site-packages/pacvis/static/MaterialIcons-Regular.svg
usr/lib/python3.10/site-packages/pacvis/static/MaterialIcons-Regular.ttf
usr/lib/python3.10/site-packages/pacvis/static/MaterialIcons-Regular.woff
usr/lib/python3.10/site-packages/pacvis/static/MaterialIcons-Regular.woff2
usr/lib/python3.10/site-packages/pacvis/static/animate.css
usr/lib/python3.10/site-packages/pacvis/static/favicon-16x16.png
usr/lib/python3.10/site-packages/pacvis/static/favicon-32x32.png
usr/lib/python3.10/site-packages/pacvis/static/material.deep_purple-amber.min.css
usr/lib/python3.10/site-packages/pacvis/static/material.min.js
usr/lib/python3.10/site-packages/pacvis/static/material.min.js.map
usr/lib/python3.10/site-packages/pacvis/static/pacvis.css
usr/lib/python3.10/site-packages/pacvis/static/pacvis.js
usr/lib/python3.10/site-packages/pacvis/static/vis-network.min.js
usr/lib/python3.10/site-packages/pacvis/static/vis.min.css
usr/lib/python3.10/site-packages/pacvis/templates/
usr/lib/python3.10/site-packages/pacvis/templates/index.template.html
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pacvis@.service
