# files/portainer-bin-1:2.13.1-1/files
%FILES%
usr/
usr/bin/
usr/bin/portainer
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/portainer.service
usr/share/
usr/share/applications/
usr/share/applications/portainer.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/portainer.png
usr/share/licenses/
usr/share/licenses/portainer-bin/
usr/share/licenses/portainer-bin/LICENSE
usr/share/portainer/
usr/share/portainer/docker-compose
usr/share/portainer/public/
usr/share/portainer/public/060b2710bdbbe3dfe48b.svg
usr/share/portainer/public/112a479c093f4729251d.png
usr/share/portainer/public/19c2b168a8d85994e14e.png
usr/share/portainer/public/29d4ee6d4a5c786588a7.svg
usr/share/portainer/public/2dcfc527d067d4ae3424.png
usr/share/portainer/public/2f517e09eb2ca6650ff5.svg
usr/share/portainer/public/347aea0738ab6bb73a56.png
usr/share/portainer/public/47026671f2ebb0d50945.png
usr/share/portainer/public/486db7cd1d4292a0d349.svg
usr/share/portainer/public/51906a7e4bd0e1d9e771.png
usr/share/portainer/public/63a301f0574f1a696ce6.png
usr/share/portainer/public/6d50eaeb9f128c130ed9.xml
usr/share/portainer/public/7a8b4f130182d19a2d7c.svg
usr/share/portainer/public/7ee8aae1b407ce0e809b.svg
usr/share/portainer/public/9d07a3ace464cfd6e27c.png
usr/share/portainer/public/c2c179f41199d25e0055.png
usr/share/portainer/public/d3a18c5002453f58c8b9.png
usr/share/portainer/public/index.html
usr/share/portainer/public/locales/
usr/share/portainer/public/locales/en/
usr/share/portainer/public/locales/en/translation.json
usr/share/portainer/public/main.a509b03c15a00f45e4d3.js
usr/share/portainer/public/main.a509b03c15a00f45e4d3.js.map
usr/share/portainer/public/main.e21146ae266f83043800.css
usr/share/portainer/public/main.e21146ae266f83043800.css.map
usr/share/portainer/public/runtime.c728d8ca1e81f0cfb05f.js
usr/share/portainer/public/runtime.c728d8ca1e81f0cfb05f.js.map
usr/share/portainer/public/timeout.html
usr/share/portainer/public/vendor.39c7b404457587ea6a6a.css
usr/share/portainer/public/vendor.39c7b404457587ea6a6a.css.map
usr/share/portainer/public/vendor.98096e683928e8563fa3.js
usr/share/portainer/public/vendor.98096e683928e8563fa3.js.LICENSE.txt
usr/share/portainer/public/vendor.98096e683928e8563fa3.js.map
