# files/matrix-registration-0.9.1-2/files
%FILES%
etc/
etc/matrix-registration/
etc/matrix-registration/alembic.ini
etc/matrix-registration/alembic/
etc/matrix-registration/alembic/env.py
etc/matrix-registration/alembic/versions/
etc/matrix-registration/alembic/versions/140a25d5f185_create_tokens_table.py
etc/matrix-registration/config.sample.yaml
usr/
usr/bin/
usr/bin/matrix-registration
usr/lib/
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/matrix_registration-0.9.0-py3.9.egg-info/
usr/lib/python3.9/site-packages/matrix_registration-0.9.0-py3.9.egg-info/PKG-INFO
usr/lib/python3.9/site-packages/matrix_registration-0.9.0-py3.9.egg-info/SOURCES.txt
usr/lib/python3.9/site-packages/matrix_registration-0.9.0-py3.9.egg-info/dependency_links.txt
usr/lib/python3.9/site-packages/matrix_registration-0.9.0-py3.9.egg-info/entry_points.txt
usr/lib/python3.9/site-packages/matrix_registration-0.9.0-py3.9.egg-info/requires.txt
usr/lib/python3.9/site-packages/matrix_registration-0.9.0-py3.9.egg-info/top_level.txt
usr/lib/python3.9/site-packages/matrix_registration/
usr/lib/python3.9/site-packages/matrix_registration/__init__.py
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/__init__.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/api.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/api.cpython-39.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/app.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/app.cpython-39.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/config.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/config.cpython-39.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/constants.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/constants.cpython-39.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/matrix_api.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/matrix_api.cpython-39.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/tokens.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/tokens.cpython-39.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/translation.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/matrix_registration/__pycache__/translation.cpython-39.pyc
usr/lib/python3.9/site-packages/matrix_registration/api.py
usr/lib/python3.9/site-packages/matrix_registration/app.py
usr/lib/python3.9/site-packages/matrix_registration/config.py
usr/lib/python3.9/site-packages/matrix_registration/constants.py
usr/lib/python3.9/site-packages/matrix_registration/matrix_api.py
usr/lib/python3.9/site-packages/matrix_registration/static/
usr/lib/python3.9/site-packages/matrix_registration/static/css/
usr/lib/python3.9/site-packages/matrix_registration/static/css/style.css
usr/lib/python3.9/site-packages/matrix_registration/static/fonts/
usr/lib/python3.9/site-packages/matrix_registration/static/fonts/Nunito_400_Nunito_700.woff2
usr/lib/python3.9/site-packages/matrix_registration/static/images/
usr/lib/python3.9/site-packages/matrix_registration/static/images/element-logo.png
usr/lib/python3.9/site-packages/matrix_registration/static/images/favicon.ico
usr/lib/python3.9/site-packages/matrix_registration/static/images/icon.png
usr/lib/python3.9/site-packages/matrix_registration/static/images/icon32x32.png
usr/lib/python3.9/site-packages/matrix_registration/static/images/riot.png
usr/lib/python3.9/site-packages/matrix_registration/static/images/tile.png
usr/lib/python3.9/site-packages/matrix_registration/static/images/valley.jpg
usr/lib/python3.9/site-packages/matrix_registration/templates/
usr/lib/python3.9/site-packages/matrix_registration/templates/register.html
usr/lib/python3.9/site-packages/matrix_registration/tokens.py
usr/lib/python3.9/site-packages/matrix_registration/translation.py
usr/lib/python3.9/site-packages/matrix_registration/translations/
usr/lib/python3.9/site-packages/matrix_registration/translations/messages.de.yaml
usr/lib/python3.9/site-packages/matrix_registration/translations/messages.en.yaml
usr/lib/python3.9/site-packages/matrix_registration/translations/messages.pt_BR.yaml
usr/lib/python3.9/site-packages/matrix_registration/translations/messages.zh_Hans.yaml
usr/lib/python3.9/site-packages/matrix_registration/wordlist.txt
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/matrix-registration.service
