# files/nextdns-unprivileged-1.37.11-1/files
%FILES%
usr/
usr/bin/
usr/bin/nextdns
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nextdns.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/nextdns.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/nextdns.conf
