# files/ktweak-git-1.0.293f748-1/files
%FILES%
etc/
etc/systemd/
etc/systemd/system/
etc/systemd/system/ktweak.service
usr/
usr/bin/
usr/bin/ktweak
