# files/snapper-support-0.1-3/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/grub-btrfs-snapper.path
usr/lib/systemd/system/grub-btrfs-snapper.service
