# files/ufw-git-0.36.r86.g506546a-1/files
%FILES%
etc/
etc/default/
etc/default/ufw
etc/ufw/
etc/ufw/after.init
etc/ufw/after.rules
etc/ufw/after6.rules
etc/ufw/applications.d/
etc/ufw/applications.d/ufw-bittorent
etc/ufw/applications.d/ufw-chat
etc/ufw/applications.d/ufw-directoryserver
etc/ufw/applications.d/ufw-dnsserver
etc/ufw/applications.d/ufw-fileserver
etc/ufw/applications.d/ufw-loginserver
etc/ufw/applications.d/ufw-mailserver
etc/ufw/applications.d/ufw-printserver
etc/ufw/applications.d/ufw-proxyserver
etc/ufw/applications.d/ufw-webserver
etc/ufw/before.init
etc/ufw/before.rules
etc/ufw/before6.rules
etc/ufw/sysctl.conf
etc/ufw/ufw.conf
etc/ufw/user.rules
etc/ufw/user6.rules
usr/
usr/bin/
usr/bin/ufw
usr/lib/
usr/lib/python3.10/
usr/lib/python3.10/site-packages/
usr/lib/python3.10/site-packages/ufw-0.37-py3.10.egg-info
usr/lib/python3.10/site-packages/ufw/
usr/lib/python3.10/site-packages/ufw/__init__.py
usr/lib/python3.10/site-packages/ufw/__pycache__/
usr/lib/python3.10/site-packages/ufw/__pycache__/__init__.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/__init__.cpython-310.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/applications.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/applications.cpython-310.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/backend.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/backend.cpython-310.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/backend_iptables.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/backend_iptables.cpython-310.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/common.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/common.cpython-310.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/frontend.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/frontend.cpython-310.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/parser.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/parser.cpython-310.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/util.cpython-310.opt-1.pyc
usr/lib/python3.10/site-packages/ufw/__pycache__/util.cpython-310.pyc
usr/lib/python3.10/site-packages/ufw/applications.py
usr/lib/python3.10/site-packages/ufw/backend.py
usr/lib/python3.10/site-packages/ufw/backend_iptables.py
usr/lib/python3.10/site-packages/ufw/common.py
usr/lib/python3.10/site-packages/ufw/frontend.py
usr/lib/python3.10/site-packages/ufw/parser.py
usr/lib/python3.10/site-packages/ufw/util.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ufw.service
usr/lib/ufw/
usr/lib/ufw/ufw-init
usr/lib/ufw/ufw-init-functions
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/ufw-framework.8.gz
usr/share/man/man8/ufw.8.gz
usr/share/ufw/
usr/share/ufw/iptables/
usr/share/ufw/iptables/after.rules
usr/share/ufw/iptables/after6.rules
usr/share/ufw/iptables/before.rules
usr/share/ufw/iptables/before6.rules
usr/share/ufw/iptables/user.rules
usr/share/ufw/iptables/user6.rules
usr/share/ufw/messages/
usr/share/ufw/messages/.gitignore
