# files/xow-0.5-2/files
%FILES%
etc/
etc/modprobe.d/
etc/modprobe.d/xow-blacklist.conf
etc/modules-load.d/
etc/modules-load.d/xow-uinput.conf
etc/udev/
etc/udev/rules.d/
etc/udev/rules.d/99-xow.rules
usr/
usr/bin/
usr/bin/xow
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/xow.service
