# files/openrgb-0.6-3/files
%FILES%
usr/
usr/bin/
usr/bin/openrgb
usr/lib/
usr/lib/modules-load.d/
usr/lib/modules-load.d/openrgb.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/openrgb.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-openrgb.rules
usr/share/
usr/share/applications/
usr/share/applications/OpenRGB.desktop
usr/share/pixmaps/
usr/share/pixmaps/OpenRGB.png
