# files/adguardhome-1:0.107.6-3/files
%FILES%
etc/
etc/default/
etc/default/adguardhome
usr/
usr/bin/
usr/bin/adguardhome
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/adguardhome.service
