# files/rustdesk-git-1.1.9.r456.g9098619-1/files
%FILES%
usr/
usr/bin/
usr/bin/rustdesk
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rustdesk.service
usr/share/
usr/share/applications/
usr/share/applications/rustdesk.desktop
usr/share/pixmaps/
usr/share/pixmaps/rustdesk.svg
usr/share/rustdesk/
usr/share/rustdesk/libsciter-gtk.so
usr/share/rustdesk/rustdesk
usr/share/rustdesk/src/
usr/share/rustdesk/src/cli.rs
usr/share/rustdesk/src/client.rs
usr/share/rustdesk/src/clipboard_file.rs
usr/share/rustdesk/src/common.rs
usr/share/rustdesk/src/ipc.rs
usr/share/rustdesk/src/lang.rs
usr/share/rustdesk/src/lang/
usr/share/rustdesk/src/lang/cn.rs
usr/share/rustdesk/src/lang/de.rs
usr/share/rustdesk/src/lang/en.rs
usr/share/rustdesk/src/lang/eo.rs
usr/share/rustdesk/src/lang/fr.rs
usr/share/rustdesk/src/lang/id.rs
usr/share/rustdesk/src/lang/it.rs
usr/share/rustdesk/src/lang/ptbr.rs
usr/share/rustdesk/src/lang/ru.rs
usr/share/rustdesk/src/lang/template.rs
usr/share/rustdesk/src/lang/tw.rs
usr/share/rustdesk/src/lib.rs
usr/share/rustdesk/src/main.rs
usr/share/rustdesk/src/platform/
usr/share/rustdesk/src/platform/linux.rs
usr/share/rustdesk/src/platform/macos.rs
usr/share/rustdesk/src/platform/mod.rs
usr/share/rustdesk/src/platform/privileges_scripts/
usr/share/rustdesk/src/platform/privileges_scripts/com.carriez.RustDesk_server.plist
usr/share/rustdesk/src/platform/privileges_scripts/com.carriez.RustDesk_service.plist
usr/share/rustdesk/src/platform/privileges_scripts/install.scpt
usr/share/rustdesk/src/platform/privileges_scripts/uninstall.scpt
usr/share/rustdesk/src/platform/windows.rs
usr/share/rustdesk/src/port_forward.rs
usr/share/rustdesk/src/rendezvous_mediator.rs
usr/share/rustdesk/src/server.rs
usr/share/rustdesk/src/server/
usr/share/rustdesk/src/server/audio_service.rs
usr/share/rustdesk/src/server/clipboard_service.rs
usr/share/rustdesk/src/server/connection.rs
usr/share/rustdesk/src/server/input_service.rs
usr/share/rustdesk/src/server/service.rs
usr/share/rustdesk/src/server/video_service.rs
usr/share/rustdesk/src/tray-icon.ico
usr/share/rustdesk/src/ui.rs
usr/share/rustdesk/src/ui/
usr/share/rustdesk/src/ui/ab.tis
usr/share/rustdesk/src/ui/chatbox.html
usr/share/rustdesk/src/ui/cm.css
usr/share/rustdesk/src/ui/cm.html
usr/share/rustdesk/src/ui/cm.rs
usr/share/rustdesk/src/ui/cm.tis
usr/share/rustdesk/src/ui/common.css
usr/share/rustdesk/src/ui/common.tis
usr/share/rustdesk/src/ui/file_transfer.css
usr/share/rustdesk/src/ui/file_transfer.tis
usr/share/rustdesk/src/ui/grid.tis
usr/share/rustdesk/src/ui/header.css
usr/share/rustdesk/src/ui/header.tis
usr/share/rustdesk/src/ui/index.css
usr/share/rustdesk/src/ui/index.html
usr/share/rustdesk/src/ui/index.tis
usr/share/rustdesk/src/ui/install.html
usr/share/rustdesk/src/ui/install.tis
usr/share/rustdesk/src/ui/macos.rs
usr/share/rustdesk/src/ui/msgbox.tis
usr/share/rustdesk/src/ui/port_forward.tis
usr/share/rustdesk/src/ui/remote.css
usr/share/rustdesk/src/ui/remote.html
usr/share/rustdesk/src/ui/remote.rs
usr/share/rustdesk/src/ui/remote.tis
usr/share/rustdesk/src/version.rs
usr/share/rustdesk/src/windows.cc
