# files/timeshift-support-0.1-3/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/grub-btrfs-timeshift.path
usr/lib/systemd/system/grub-btrfs-timeshift.service
usr/share/
usr/share/libalpm/
usr/share/libalpm/hooks/
usr/share/libalpm/hooks/timeshift-autosnap-fix.hook
