# files/arch-audit-git-0.1.20.r0.gb9f7220-1/files
%FILES%
usr/
usr/bin/
usr/bin/arch-audit
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/arch-audit.service
usr/lib/systemd/system/arch-audit.timer
usr/share/
usr/share/arch-audit-git/
usr/share/arch-audit-git/arch-audit.hook
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/arch-audit
usr/share/doc/
usr/share/doc/arch-audit-git/
usr/share/doc/arch-audit-git/README.md
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/arch-audit.fish
usr/share/licenses/
usr/share/licenses/arch-audit-git/
usr/share/licenses/arch-audit-git/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/arch-audit.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_arch-audit
