# files/virt-manager-meta-4-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/libvirt.service.d/
usr/lib/systemd/system/libvirt.service.d/firewall.conf
