# files/cassandra-4.0.0-1/files
%FILES%
etc/
etc/cassandra/
etc/cassandra/README.txt
etc/cassandra/cassandra-env.sh
etc/cassandra/cassandra-jaas.config
etc/cassandra/cassandra-rackdc.properties
etc/cassandra/cassandra-topology.properties
etc/cassandra/cassandra.yaml
etc/cassandra/commitlog_archiving.properties
etc/cassandra/cqlshrc.sample
etc/cassandra/hotspot_compiler
etc/cassandra/jvm-clients.options
etc/cassandra/jvm-server.options
etc/cassandra/jvm11-clients.options
etc/cassandra/jvm11-server.options
etc/cassandra/jvm8-clients.options
etc/cassandra/jvm8-server.options
etc/cassandra/logback-tools.xml
etc/cassandra/logback.xml
etc/cassandra/metrics-reporter-config-sample.yaml
etc/cassandra/triggers/
etc/cassandra/triggers/README.txt
usr/
usr/bin/
usr/bin/cassandra
usr/bin/cqlsh
usr/bin/debug-cql
usr/bin/nodetool
usr/bin/sstableloader
usr/bin/sstablescrub
usr/bin/sstableupgrade
usr/bin/sstableutil
usr/bin/sstableverify
usr/bin/stop-server
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/cassandra.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/cassandra.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/cassandra.conf
usr/share/
usr/share/cassandra/
usr/share/cassandra/bin/
usr/share/cassandra/bin/cassandra
usr/share/cassandra/bin/cqlsh
usr/share/cassandra/bin/cqlsh.py
usr/share/cassandra/bin/debug-cql
usr/share/cassandra/bin/nodetool
usr/share/cassandra/bin/sstableloader
usr/share/cassandra/bin/sstablescrub
usr/share/cassandra/bin/sstableupgrade
usr/share/cassandra/bin/sstableutil
usr/share/cassandra/bin/sstableverify
usr/share/cassandra/bin/stop-server
usr/share/cassandra/cassandra.in.sh
usr/share/cassandra/conf
usr/share/cassandra/lib
usr/share/cassandra/pylib/
usr/share/cassandra/pylib/Dockerfile.ubuntu.py2
usr/share/cassandra/pylib/Dockerfile.ubuntu.py3
usr/share/cassandra/pylib/Dockerfile.ubuntu.py37
usr/share/cassandra/pylib/Dockerfile.ubuntu.py38
usr/share/cassandra/pylib/README.asc
usr/share/cassandra/pylib/cassandra-cqlsh-tests.sh
usr/share/cassandra/pylib/cqlshlib/
usr/share/cassandra/pylib/cqlshlib/__init__.py
usr/share/cassandra/pylib/cqlshlib/copyutil.py
usr/share/cassandra/pylib/cqlshlib/cql3handling.py
usr/share/cassandra/pylib/cqlshlib/cqlhandling.py
usr/share/cassandra/pylib/cqlshlib/cqlshhandling.py
usr/share/cassandra/pylib/cqlshlib/displaying.py
usr/share/cassandra/pylib/cqlshlib/formatting.py
usr/share/cassandra/pylib/cqlshlib/helptopics.py
usr/share/cassandra/pylib/cqlshlib/pylexotron.py
usr/share/cassandra/pylib/cqlshlib/saferscanner.py
usr/share/cassandra/pylib/cqlshlib/setup.cfg
usr/share/cassandra/pylib/cqlshlib/sslhandling.py
usr/share/cassandra/pylib/cqlshlib/test/
usr/share/cassandra/pylib/cqlshlib/test/__init__.py
usr/share/cassandra/pylib/cqlshlib/test/ansi_colors.py
usr/share/cassandra/pylib/cqlshlib/test/basecase.py
usr/share/cassandra/pylib/cqlshlib/test/cassconnect.py
usr/share/cassandra/pylib/cqlshlib/test/config/
usr/share/cassandra/pylib/cqlshlib/test/config/sslhandling.config
usr/share/cassandra/pylib/cqlshlib/test/config/sslhandling_invalid.config
usr/share/cassandra/pylib/cqlshlib/test/run_cqlsh.py
usr/share/cassandra/pylib/cqlshlib/test/test_constants.py
usr/share/cassandra/pylib/cqlshlib/test/test_copyutil.py
usr/share/cassandra/pylib/cqlshlib/test/test_cql_parsing.py
usr/share/cassandra/pylib/cqlshlib/test/test_cqlsh_completion.py
usr/share/cassandra/pylib/cqlshlib/test/test_cqlsh_output.py
usr/share/cassandra/pylib/cqlshlib/test/test_keyspace_init.cql
usr/share/cassandra/pylib/cqlshlib/test/test_sslhandling.py
usr/share/cassandra/pylib/cqlshlib/test/test_unicode.py
usr/share/cassandra/pylib/cqlshlib/test/winpty.py
usr/share/cassandra/pylib/cqlshlib/tracing.py
usr/share/cassandra/pylib/cqlshlib/util.py
usr/share/cassandra/pylib/cqlshlib/wcwidth.py
usr/share/cassandra/pylib/requirements.txt
usr/share/cassandra/pylib/setup.py
usr/share/cassandra/tools/
usr/share/cassandra/tools/bin/
usr/share/cassandra/tools/bin/auditlogviewer
usr/share/cassandra/tools/bin/cassandra-stress
usr/share/cassandra/tools/bin/cassandra-stressd
usr/share/cassandra/tools/bin/cassandra.in.sh
usr/share/cassandra/tools/bin/compaction-stress
usr/share/cassandra/tools/bin/fqltool
usr/share/cassandra/tools/bin/generatetokens
usr/share/cassandra/tools/bin/jmxtool
usr/share/cassandra/tools/bin/sstabledump
usr/share/cassandra/tools/bin/sstableexpiredblockers
usr/share/cassandra/tools/bin/sstablelevelreset
usr/share/cassandra/tools/bin/sstablemetadata
usr/share/cassandra/tools/bin/sstableofflinerelevel
usr/share/cassandra/tools/bin/sstablerepairedset
usr/share/cassandra/tools/bin/sstablesplit
usr/share/cassandra/tools/cqlstress-counter-example.yaml
usr/share/cassandra/tools/cqlstress-example.yaml
usr/share/cassandra/tools/cqlstress-insanity-example.yaml
usr/share/cassandra/tools/cqlstress-lwt-example.yaml
usr/share/cassandra/tools/lib/
usr/share/cassandra/tools/lib/fqltool.jar
usr/share/cassandra/tools/lib/stress.jar
usr/share/java/
usr/share/java/cassandra/
usr/share/java/cassandra/HdrHistogram-2.1.9.jar
usr/share/java/cassandra/ST4-4.0.8.jar
usr/share/java/cassandra/airline-0.8.jar
usr/share/java/cassandra/antlr-runtime-3.5.2.jar
usr/share/java/cassandra/apache-cassandra-4.0.0.jar
usr/share/java/cassandra/asm-7.1.jar
usr/share/java/cassandra/caffeine-2.3.5.jar
usr/share/java/cassandra/cassandra-driver-core-3.11.0-shaded.jar
usr/share/java/cassandra/cassandra-driver-internal-only-3.25.0.zip
usr/share/java/cassandra/chronicle-bytes-2.20.111.jar
usr/share/java/cassandra/chronicle-core-2.20.126.jar
usr/share/java/cassandra/chronicle-queue-5.20.123.jar
usr/share/java/cassandra/chronicle-threads-2.20.111.jar
usr/share/java/cassandra/chronicle-wire-2.20.117.jar
usr/share/java/cassandra/commons-cli-1.1.jar
usr/share/java/cassandra/commons-codec-1.9.jar
usr/share/java/cassandra/commons-lang3-3.11.jar
usr/share/java/cassandra/commons-math3-3.2.jar
usr/share/java/cassandra/concurrent-trees-2.4.0.jar
usr/share/java/cassandra/ecj-4.6.1.jar
usr/share/java/cassandra/futures-2.1.6-py2.py3-none-any.zip
usr/share/java/cassandra/geomet-0.1.0.zip
usr/share/java/cassandra/guava-27.0-jre.jar
usr/share/java/cassandra/high-scale-lib-1.0.6.jar
usr/share/java/cassandra/hppc-0.8.1.jar
usr/share/java/cassandra/j2objc-annotations-1.3.jar
usr/share/java/cassandra/jackson-annotations-2.9.10.jar
usr/share/java/cassandra/jackson-core-2.9.10.jar
usr/share/java/cassandra/jackson-databind-2.9.10.8.jar
usr/share/java/cassandra/jamm-0.3.2.jar
usr/share/java/cassandra/java-cup-runtime-11b-20160615.jar
usr/share/java/cassandra/javax.inject-1.jar
usr/share/java/cassandra/jbcrypt-0.3m.jar
usr/share/java/cassandra/jcl-over-slf4j-1.7.25.jar
usr/share/java/cassandra/jcommander-1.30.jar
usr/share/java/cassandra/jctools-core-3.1.0.jar
usr/share/java/cassandra/jflex-1.8.2.jar
usr/share/java/cassandra/jna-5.6.0.jar
usr/share/java/cassandra/json-simple-1.1.jar
usr/share/java/cassandra/jvm-attach-api-1.5.jar
usr/share/java/cassandra/log4j-over-slf4j-1.7.25.jar
usr/share/java/cassandra/logback-classic-1.2.3.jar
usr/share/java/cassandra/logback-core-1.2.3.jar
usr/share/java/cassandra/lz4-java-1.7.1.jar
usr/share/java/cassandra/metrics-core-3.1.5.jar
usr/share/java/cassandra/metrics-jvm-3.1.5.jar
usr/share/java/cassandra/metrics-logback-3.1.5.jar
usr/share/java/cassandra/mxdump-0.14.jar
usr/share/java/cassandra/netty-all-4.1.58.Final.jar
usr/share/java/cassandra/netty-tcnative-boringssl-static-2.0.36.Final.jar
usr/share/java/cassandra/ohc-core-0.5.1.jar
usr/share/java/cassandra/ohc-core-j8-0.5.1.jar
usr/share/java/cassandra/psjava-0.1.19.jar
usr/share/java/cassandra/reporter-config-base-3.0.3.jar
usr/share/java/cassandra/reporter-config3-3.0.3.jar
usr/share/java/cassandra/sigar-1.6.4.jar
usr/share/java/cassandra/sigar-bin/
usr/share/java/cassandra/sigar-bin/libsigar-amd64-freebsd-6.so
usr/share/java/cassandra/sigar-bin/libsigar-amd64-linux.so
usr/share/java/cassandra/sigar-bin/libsigar-amd64-solaris.so
usr/share/java/cassandra/sigar-bin/libsigar-ia64-hpux-11.sl
usr/share/java/cassandra/sigar-bin/libsigar-ia64-linux.so
usr/share/java/cassandra/sigar-bin/libsigar-pa-hpux-11.sl
usr/share/java/cassandra/sigar-bin/libsigar-ppc-aix-5.so
usr/share/java/cassandra/sigar-bin/libsigar-ppc-linux.so
usr/share/java/cassandra/sigar-bin/libsigar-ppc64-aix-5.so
usr/share/java/cassandra/sigar-bin/libsigar-ppc64-linux.so
usr/share/java/cassandra/sigar-bin/libsigar-ppc64le-linux.so
usr/share/java/cassandra/sigar-bin/libsigar-s390x-linux.so
usr/share/java/cassandra/sigar-bin/libsigar-sparc-solaris.so
usr/share/java/cassandra/sigar-bin/libsigar-sparc64-solaris.so
usr/share/java/cassandra/sigar-bin/libsigar-universal-macosx.dylib
usr/share/java/cassandra/sigar-bin/libsigar-universal64-macosx.dylib
usr/share/java/cassandra/sigar-bin/libsigar-x86-freebsd-5.so
usr/share/java/cassandra/sigar-bin/libsigar-x86-freebsd-6.so
usr/share/java/cassandra/sigar-bin/libsigar-x86-linux.so
usr/share/java/cassandra/sigar-bin/libsigar-x86-solaris.so
usr/share/java/cassandra/sigar-bin/sigar-amd64-winnt.dll
usr/share/java/cassandra/sigar-bin/sigar-x86-winnt.dll
usr/share/java/cassandra/sigar-bin/sigar-x86-winnt.lib
usr/share/java/cassandra/six-1.12.0-py2.py3-none-any.zip
usr/share/java/cassandra/sjk-cli-0.14.jar
usr/share/java/cassandra/sjk-core-0.14.jar
usr/share/java/cassandra/sjk-json-0.14.jar
usr/share/java/cassandra/sjk-stacktrace-0.14.jar
usr/share/java/cassandra/slf4j-api-1.7.25.jar
usr/share/java/cassandra/snakeyaml-1.26.jar
usr/share/java/cassandra/snappy-java-1.1.2.6.jar
usr/share/java/cassandra/snowball-stemmer-1.3.0.581.1.jar
usr/share/java/cassandra/stream-2.5.2.jar
usr/share/java/cassandra/zstd-jni-1.3.8-5.jar
var/
var/log/
var/log/cassandra/
