# files/qemu-guest-agent-git-14:6.1.0.r1419.ge5b2333f24-1/files
%FILES%
etc/
etc/qemu/
etc/qemu/fsfreeze-hook
usr/
usr/bin/
usr/bin/qga
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/qemu-guest-agent.service
