# files/xmm7360-pci-utils-git-r229.cf6625a-1/files
%FILES%
etc/
etc/xmm7360
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/xmm7360.service
usr/lib/xmm7360-pci-git/
usr/lib/xmm7360-pci-git/rpc/
usr/lib/xmm7360-pci-git/rpc/.gitignore
usr/lib/xmm7360-pci-git/rpc/Makefile
usr/lib/xmm7360-pci-git/rpc/__pycache__/
usr/lib/xmm7360-pci-git/rpc/__pycache__/mux.cpython-310.opt-1.pyc
usr/lib/xmm7360-pci-git/rpc/__pycache__/open_xdatachannel.cpython-310.opt-1.pyc
usr/lib/xmm7360-pci-git/rpc/__pycache__/rpc.cpython-310.opt-1.pyc
usr/lib/xmm7360-pci-git/rpc/__pycache__/rpc_call_ids.cpython-310.opt-1.pyc
usr/lib/xmm7360-pci-git/rpc/__pycache__/rpc_unsol_table.cpython-310.opt-1.pyc
usr/lib/xmm7360-pci-git/rpc/__pycache__/test_rpc.cpython-310.opt-1.pyc
usr/lib/xmm7360-pci-git/rpc/__pycache__/xm_dbus.cpython-310.opt-1.pyc
usr/lib/xmm7360-pci-git/rpc/mux.c
usr/lib/xmm7360-pci-git/rpc/mux.py
usr/lib/xmm7360-pci-git/rpc/open_xdatachannel.py
usr/lib/xmm7360-pci-git/rpc/rpc.py
usr/lib/xmm7360-pci-git/rpc/rpc_call_ids.py
usr/lib/xmm7360-pci-git/rpc/rpc_unsol_table.py
usr/lib/xmm7360-pci-git/rpc/test_rpc.py
usr/lib/xmm7360-pci-git/rpc/xm_dbus.py
usr/lib/xmm7360-pci-git/rpc/xmm7360.h
