# files/gtk3-nocsd-git-r82.512c2bd-1/files
%FILES%
etc/
etc/X11/
etc/X11/xinit/
etc/X11/xinit/xinitrc.d/
etc/X11/xinit/xinitrc.d/30-gtk3-nocsd.sh
usr/
usr/bin/
usr/bin/gtk3-nocsd
usr/lib/
usr/lib/libgtk3-nocsd.so.0
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/gtk3-nocsd
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/gtk3-nocsd.1.gz
