# files/wireproxy-1.0.5-2/files
%FILES%
usr/
usr/bin/
usr/bin/wireproxy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/wireproxy.service
usr/lib/systemd/system/wireproxy@.service
usr/lib/systemd/user/
usr/lib/systemd/user/wireproxy.service
usr/lib/systemd/user/wireproxy@.service
usr/share/
usr/share/licenses/
usr/share/licenses/wireproxy/
usr/share/licenses/wireproxy/LICENSE
