# files/anope-2.0.11-1/files
%FILES%
etc/
etc/anope/
etc/anope/botserv.example.conf
etc/anope/chanserv.example.conf
etc/anope/chanstats.example.conf
etc/anope/example.conf
etc/anope/global.example.conf
etc/anope/hostserv.example.conf
etc/anope/irc2sql.example.conf
etc/anope/memoserv.example.conf
etc/anope/modules.example.conf
etc/anope/nickserv.example.conf
etc/anope/operserv.example.conf
etc/anope/stats.standalone.example.conf
etc/tmpfiles.d/
etc/tmpfiles.d/anope.conf
usr/
usr/bin/
usr/bin/anopesmtp
usr/bin/services
usr/lib/
usr/lib/anope/
usr/lib/anope/locale/
usr/lib/anope/locale/ca_ES/
usr/lib/anope/locale/ca_ES/LC_MESSAGES/
usr/lib/anope/locale/ca_ES/LC_MESSAGES/anope.mo
usr/lib/anope/locale/de_DE/
usr/lib/anope/locale/de_DE/LC_MESSAGES/
usr/lib/anope/locale/de_DE/LC_MESSAGES/anope.mo
usr/lib/anope/locale/el_GR/
usr/lib/anope/locale/el_GR/LC_MESSAGES/
usr/lib/anope/locale/el_GR/LC_MESSAGES/anope.mo
usr/lib/anope/locale/en_US/
usr/lib/anope/locale/en_US/LC_MESSAGES/
usr/lib/anope/locale/en_US/LC_MESSAGES/anope.mo
usr/lib/anope/locale/es_ES/
usr/lib/anope/locale/es_ES/LC_MESSAGES/
usr/lib/anope/locale/es_ES/LC_MESSAGES/anope.mo
usr/lib/anope/locale/fr_FR/
usr/lib/anope/locale/fr_FR/LC_MESSAGES/
usr/lib/anope/locale/fr_FR/LC_MESSAGES/anope.mo
usr/lib/anope/locale/hu_HU/
usr/lib/anope/locale/hu_HU/LC_MESSAGES/
usr/lib/anope/locale/hu_HU/LC_MESSAGES/anope.mo
usr/lib/anope/locale/it_IT/
usr/lib/anope/locale/it_IT/LC_MESSAGES/
usr/lib/anope/locale/it_IT/LC_MESSAGES/anope.mo
usr/lib/anope/locale/nl_NL/
usr/lib/anope/locale/nl_NL/LC_MESSAGES/
usr/lib/anope/locale/nl_NL/LC_MESSAGES/anope.mo
usr/lib/anope/locale/pl_PL/
usr/lib/anope/locale/pl_PL/LC_MESSAGES/
usr/lib/anope/locale/pl_PL/LC_MESSAGES/anope.mo
usr/lib/anope/locale/pt_PT/
usr/lib/anope/locale/pt_PT/LC_MESSAGES/
usr/lib/anope/locale/pt_PT/LC_MESSAGES/anope.mo
usr/lib/anope/locale/ru_RU/
usr/lib/anope/locale/ru_RU/LC_MESSAGES/
usr/lib/anope/locale/ru_RU/LC_MESSAGES/anope.mo
usr/lib/anope/locale/tr_TR/
usr/lib/anope/locale/tr_TR/LC_MESSAGES/
usr/lib/anope/locale/tr_TR/LC_MESSAGES/anope.mo
usr/lib/anope/modules/
usr/lib/anope/modules/bahamut.so
usr/lib/anope/modules/botserv.so
usr/lib/anope/modules/bs_assign.so
usr/lib/anope/modules/bs_autoassign.so
usr/lib/anope/modules/bs_badwords.so
usr/lib/anope/modules/bs_bot.so
usr/lib/anope/modules/bs_botlist.so
usr/lib/anope/modules/bs_control.so
usr/lib/anope/modules/bs_info.so
usr/lib/anope/modules/bs_kick.so
usr/lib/anope/modules/bs_set.so
usr/lib/anope/modules/chanserv.so
usr/lib/anope/modules/charybdis.so
usr/lib/anope/modules/cs_access.so
usr/lib/anope/modules/cs_akick.so
usr/lib/anope/modules/cs_ban.so
usr/lib/anope/modules/cs_clone.so
usr/lib/anope/modules/cs_drop.so
usr/lib/anope/modules/cs_enforce.so
usr/lib/anope/modules/cs_entrymsg.so
usr/lib/anope/modules/cs_flags.so
usr/lib/anope/modules/cs_getkey.so
usr/lib/anope/modules/cs_info.so
usr/lib/anope/modules/cs_invite.so
usr/lib/anope/modules/cs_kick.so
usr/lib/anope/modules/cs_list.so
usr/lib/anope/modules/cs_log.so
usr/lib/anope/modules/cs_mode.so
usr/lib/anope/modules/cs_register.so
usr/lib/anope/modules/cs_seen.so
usr/lib/anope/modules/cs_set.so
usr/lib/anope/modules/cs_set_misc.so
usr/lib/anope/modules/cs_status.so
usr/lib/anope/modules/cs_statusupdate.so
usr/lib/anope/modules/cs_suspend.so
usr/lib/anope/modules/cs_sync.so
usr/lib/anope/modules/cs_topic.so
usr/lib/anope/modules/cs_unban.so
usr/lib/anope/modules/cs_updown.so
usr/lib/anope/modules/cs_xop.so
usr/lib/anope/modules/db_flatfile.so
usr/lib/anope/modules/db_old.so
usr/lib/anope/modules/db_redis.so
usr/lib/anope/modules/db_sql.so
usr/lib/anope/modules/db_sql_live.so
usr/lib/anope/modules/enc_bcrypt.so
usr/lib/anope/modules/enc_md5.so
usr/lib/anope/modules/enc_none.so
usr/lib/anope/modules/enc_old.so
usr/lib/anope/modules/enc_sha1.so
usr/lib/anope/modules/enc_sha256.so
usr/lib/anope/modules/fantasy.so
usr/lib/anope/modules/gl_global.so
usr/lib/anope/modules/global.so
usr/lib/anope/modules/greet.so
usr/lib/anope/modules/help.so
usr/lib/anope/modules/hostserv.so
usr/lib/anope/modules/hs_del.so
usr/lib/anope/modules/hs_group.so
usr/lib/anope/modules/hs_list.so
usr/lib/anope/modules/hs_off.so
usr/lib/anope/modules/hs_on.so
usr/lib/anope/modules/hs_request.so
usr/lib/anope/modules/hs_set.so
usr/lib/anope/modules/hybrid.so
usr/lib/anope/modules/inspircd12.so
usr/lib/anope/modules/inspircd20.so
usr/lib/anope/modules/inspircd3.so
usr/lib/anope/modules/m_dns.so
usr/lib/anope/modules/m_dnsbl.so
usr/lib/anope/modules/m_helpchan.so
usr/lib/anope/modules/m_httpd.so
usr/lib/anope/modules/m_proxyscan.so
usr/lib/anope/modules/m_redis.so
usr/lib/anope/modules/m_rewrite.so
usr/lib/anope/modules/m_sasl.so
usr/lib/anope/modules/m_xmlrpc.so
usr/lib/anope/modules/m_xmlrpc_main.so
usr/lib/anope/modules/memoserv.so
usr/lib/anope/modules/ms_cancel.so
usr/lib/anope/modules/ms_check.so
usr/lib/anope/modules/ms_del.so
usr/lib/anope/modules/ms_ignore.so
usr/lib/anope/modules/ms_info.so
usr/lib/anope/modules/ms_list.so
usr/lib/anope/modules/ms_read.so
usr/lib/anope/modules/ms_rsend.so
usr/lib/anope/modules/ms_send.so
usr/lib/anope/modules/ms_sendall.so
usr/lib/anope/modules/ms_set.so
usr/lib/anope/modules/ms_staff.so
usr/lib/anope/modules/ngircd.so
usr/lib/anope/modules/nickserv.so
usr/lib/anope/modules/ns_access.so
usr/lib/anope/modules/ns_ajoin.so
usr/lib/anope/modules/ns_alist.so
usr/lib/anope/modules/ns_cert.so
usr/lib/anope/modules/ns_drop.so
usr/lib/anope/modules/ns_getemail.so
usr/lib/anope/modules/ns_getpass.so
usr/lib/anope/modules/ns_group.so
usr/lib/anope/modules/ns_identify.so
usr/lib/anope/modules/ns_info.so
usr/lib/anope/modules/ns_list.so
usr/lib/anope/modules/ns_logout.so
usr/lib/anope/modules/ns_maxemail.so
usr/lib/anope/modules/ns_recover.so
usr/lib/anope/modules/ns_register.so
usr/lib/anope/modules/ns_resetpass.so
usr/lib/anope/modules/ns_set.so
usr/lib/anope/modules/ns_set_misc.so
usr/lib/anope/modules/ns_status.so
usr/lib/anope/modules/ns_suspend.so
usr/lib/anope/modules/ns_update.so
usr/lib/anope/modules/operserv.so
usr/lib/anope/modules/os_akill.so
usr/lib/anope/modules/os_chankill.so
usr/lib/anope/modules/os_config.so
usr/lib/anope/modules/os_defcon.so
usr/lib/anope/modules/os_dns.so
usr/lib/anope/modules/os_forbid.so
usr/lib/anope/modules/os_ignore.so
usr/lib/anope/modules/os_info.so
usr/lib/anope/modules/os_jupe.so
usr/lib/anope/modules/os_kick.so
usr/lib/anope/modules/os_kill.so
usr/lib/anope/modules/os_list.so
usr/lib/anope/modules/os_login.so
usr/lib/anope/modules/os_logsearch.so
usr/lib/anope/modules/os_mode.so
usr/lib/anope/modules/os_modinfo.so
usr/lib/anope/modules/os_module.so
usr/lib/anope/modules/os_news.so
usr/lib/anope/modules/os_noop.so
usr/lib/anope/modules/os_oline.so
usr/lib/anope/modules/os_oper.so
usr/lib/anope/modules/os_reload.so
usr/lib/anope/modules/os_session.so
usr/lib/anope/modules/os_set.so
usr/lib/anope/modules/os_shutdown.so
usr/lib/anope/modules/os_stats.so
usr/lib/anope/modules/os_svs.so
usr/lib/anope/modules/os_sxline.so
usr/lib/anope/modules/os_update.so
usr/lib/anope/modules/plexus.so
usr/lib/anope/modules/ratbox.so
usr/lib/anope/modules/unreal.so
usr/lib/anope/modules/unreal4.so
usr/lib/anope/modules/webcpanel.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/anope.service
var/
var/lib/
var/lib/anope/
var/lib/anope/backups/
var/lib/anope/modules/
var/lib/anope/modules/webcpanel/
var/lib/anope/modules/webcpanel/templates/
var/lib/anope/modules/webcpanel/templates/default/
var/lib/anope/modules/webcpanel/templates/default/chanserv/
var/lib/anope/modules/webcpanel/templates/default/chanserv/access.html
var/lib/anope/modules/webcpanel/templates/default/chanserv/akick.html
var/lib/anope/modules/webcpanel/templates/default/chanserv/chanlist.html
var/lib/anope/modules/webcpanel/templates/default/chanserv/drop.html
var/lib/anope/modules/webcpanel/templates/default/chanserv/main.html
var/lib/anope/modules/webcpanel/templates/default/chanserv/modes.html
var/lib/anope/modules/webcpanel/templates/default/chanserv/set.html
var/lib/anope/modules/webcpanel/templates/default/confirm.html
var/lib/anope/modules/webcpanel/templates/default/cubes.png
var/lib/anope/modules/webcpanel/templates/default/favicon.ico
var/lib/anope/modules/webcpanel/templates/default/footer.html
var/lib/anope/modules/webcpanel/templates/default/header.html
var/lib/anope/modules/webcpanel/templates/default/hostserv/
var/lib/anope/modules/webcpanel/templates/default/hostserv/request.html
var/lib/anope/modules/webcpanel/templates/default/login.html
var/lib/anope/modules/webcpanel/templates/default/logo.png
var/lib/anope/modules/webcpanel/templates/default/memoserv/
var/lib/anope/modules/webcpanel/templates/default/memoserv/memos.html
var/lib/anope/modules/webcpanel/templates/default/nickserv/
var/lib/anope/modules/webcpanel/templates/default/nickserv/access.html
var/lib/anope/modules/webcpanel/templates/default/nickserv/alist.html
var/lib/anope/modules/webcpanel/templates/default/nickserv/cert.html
var/lib/anope/modules/webcpanel/templates/default/nickserv/info.html
var/lib/anope/modules/webcpanel/templates/default/operserv/
var/lib/anope/modules/webcpanel/templates/default/operserv/akill.html
var/lib/anope/modules/webcpanel/templates/default/register.html
var/lib/anope/modules/webcpanel/templates/default/style.css
var/lib/anope/runtime/
var/log/
var/log/anope/
