# files/logmein-hamachi-2.1.0.203-1/files
%FILES%
opt/
opt/logmein-hamachi/
opt/logmein-hamachi/CHANGES
opt/logmein-hamachi/LICENSE
opt/logmein-hamachi/README
opt/logmein-hamachi/bin/
opt/logmein-hamachi/bin/dnsdown
opt/logmein-hamachi/bin/dnsup
opt/logmein-hamachi/bin/hamachid
usr/
usr/bin/
usr/bin/hamachi
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/logmein-hamachi.service
usr/share/
usr/share/licenses/
usr/share/licenses/logmein-hamachi/
usr/share/licenses/logmein-hamachi/LICENSE
