# files/clamav-unofficial-sigs-7.2.5-1/files
%FILES%
etc/
etc/clamav-unofficial-sigs/
etc/clamav-unofficial-sigs/master.conf
etc/clamav-unofficial-sigs/os.conf
etc/clamav-unofficial-sigs/user.conf
etc/logrotate.d/
etc/logrotate.d/clamav-unofficial-sigs
usr/
usr/bin/
usr/bin/clamav-unofficial-sigs.sh
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/clamav-unofficial-sigs.service
usr/lib/systemd/system/clamav-unofficial-sigs.timer
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/clamav-unofficial-sigs.conf
usr/share/
usr/share/licenses/
usr/share/licenses/clamav-unofficial-sigs/
usr/share/licenses/clamav-unofficial-sigs/LICENSE
