# files/cloudflared-2021.9.1-1/files
%FILES%
etc/
etc/cloudflared/
etc/cloudflared/config.yml
usr/
usr/bin/
usr/bin/cloudflared
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/cloudflared.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/cloudflared.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/cloudflared.conf
usr/share/
usr/share/licenses/
usr/share/licenses/cloudflared/
usr/share/licenses/cloudflared/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/cloudflared.1.gz
