# files/sddm-git-0.19.0.63.g83e1ae5-1/files
%FILES%
etc/
etc/pam.d/
etc/pam.d/sddm
etc/pam.d/sddm-autologin
etc/pam.d/sddm-greeter
usr/
usr/bin/
usr/bin/sddm
usr/bin/sddm-greeter
usr/lib/
usr/lib/qt/
usr/lib/qt/qml/
usr/lib/qt/qml/SddmComponents/
usr/lib/qt/qml/SddmComponents/Background.qml
usr/lib/qt/qml/SddmComponents/Button.qml
usr/lib/qt/qml/SddmComponents/Clock.qml
usr/lib/qt/qml/SddmComponents/ComboBox.qml
usr/lib/qt/qml/SddmComponents/ImageButton.qml
usr/lib/qt/qml/SddmComponents/LayoutBox.qml
usr/lib/qt/qml/SddmComponents/Menu.qml
usr/lib/qt/qml/SddmComponents/PasswordBox.qml
usr/lib/qt/qml/SddmComponents/PictureBox.qml
usr/lib/qt/qml/SddmComponents/TextBox.qml
usr/lib/qt/qml/SddmComponents/TextConstants.qml
usr/lib/qt/qml/SddmComponents/qmldir
usr/lib/qt/qml/SddmComponents/warning.png
usr/lib/sddm/
usr/lib/sddm/sddm-helper
usr/lib/sddm/sddm.conf.d/
usr/lib/sddm/sddm.conf.d/default.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sddm.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/sddm.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/sddm.conf
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/sddm_org.freedesktop.DisplayManager.conf
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/sddm-greeter.1.gz
usr/share/man/man1/sddm.1.gz
usr/share/man/man5/
usr/share/man/man5/sddm-state.conf.5.gz
usr/share/man/man5/sddm.conf.5.gz
usr/share/sddm/
usr/share/sddm/faces/
usr/share/sddm/faces/.face.icon
usr/share/sddm/faces/root.face.icon
usr/share/sddm/flags/
usr/share/sddm/flags/ae.png
usr/share/sddm/flags/am.png
usr/share/sddm/flags/ar.png
usr/share/sddm/flags/at.png
usr/share/sddm/flags/az.png
usr/share/sddm/flags/bd.png
usr/share/sddm/flags/be.png
usr/share/sddm/flags/bg.png
usr/share/sddm/flags/bh.png
usr/share/sddm/flags/br.png
usr/share/sddm/flags/by.png
usr/share/sddm/flags/ca.png
usr/share/sddm/flags/ch.png
usr/share/sddm/flags/cu.png
usr/share/sddm/flags/cz.png
usr/share/sddm/flags/de.png
usr/share/sddm/flags/dj.png
usr/share/sddm/flags/dk.png
usr/share/sddm/flags/dz.png
usr/share/sddm/flags/ee.png
usr/share/sddm/flags/eg.png
usr/share/sddm/flags/es.png
usr/share/sddm/flags/eu.png
usr/share/sddm/flags/fi.png
usr/share/sddm/flags/fr.png
usr/share/sddm/flags/gb.png
usr/share/sddm/flags/ge.png
usr/share/sddm/flags/gr.png
usr/share/sddm/flags/hr.png
usr/share/sddm/flags/hu.png
usr/share/sddm/flags/il.png
usr/share/sddm/flags/in.png
usr/share/sddm/flags/iq.png
usr/share/sddm/flags/is.png
usr/share/sddm/flags/it.png
usr/share/sddm/flags/jo.png
usr/share/sddm/flags/jp.png
usr/share/sddm/flags/km.png
usr/share/sddm/flags/kr.png
usr/share/sddm/flags/kw.png
usr/share/sddm/flags/la.png
usr/share/sddm/flags/lb.png
usr/share/sddm/flags/lt.png
usr/share/sddm/flags/lv.png
usr/share/sddm/flags/ly.png
usr/share/sddm/flags/ma.png
usr/share/sddm/flags/mk.png
usr/share/sddm/flags/mn.png
usr/share/sddm/flags/mx.png
usr/share/sddm/flags/nl.png
usr/share/sddm/flags/no.png
usr/share/sddm/flags/om.png
usr/share/sddm/flags/pl.png
usr/share/sddm/flags/ps.png
usr/share/sddm/flags/pt.png
usr/share/sddm/flags/qa.png
usr/share/sddm/flags/qc.png
usr/share/sddm/flags/ro.png
usr/share/sddm/flags/ru.png
usr/share/sddm/flags/sa.png
usr/share/sddm/flags/sd.png
usr/share/sddm/flags/se.png
usr/share/sddm/flags/si.png
usr/share/sddm/flags/sk.png
usr/share/sddm/flags/so.png
usr/share/sddm/flags/sr.png
usr/share/sddm/flags/sy.png
usr/share/sddm/flags/th.png
usr/share/sddm/flags/tn.png
usr/share/sddm/flags/tr.png
usr/share/sddm/flags/ua.png
usr/share/sddm/flags/uk.png
usr/share/sddm/flags/un.png
usr/share/sddm/flags/us.png
usr/share/sddm/flags/uy.png
usr/share/sddm/flags/vn.png
usr/share/sddm/flags/ye.png
usr/share/sddm/flags/yu.png
usr/share/sddm/flags/zz.png
usr/share/sddm/scripts/
usr/share/sddm/scripts/Xsession
usr/share/sddm/scripts/Xsetup
usr/share/sddm/scripts/Xstop
usr/share/sddm/scripts/wayland-session
usr/share/sddm/themes/
usr/share/sddm/themes/elarun/
usr/share/sddm/themes/elarun/Main.qml
usr/share/sddm/themes/elarun/README
usr/share/sddm/themes/elarun/angle-down.png
usr/share/sddm/themes/elarun/elarun.jpg
usr/share/sddm/themes/elarun/images/
usr/share/sddm/themes/elarun/images/background.png
usr/share/sddm/themes/elarun/images/lock.png
usr/share/sddm/themes/elarun/images/login_active.png
usr/share/sddm/themes/elarun/images/login_normal.png
usr/share/sddm/themes/elarun/images/rectangle.png
usr/share/sddm/themes/elarun/images/rectangle_overlay.png
usr/share/sddm/themes/elarun/images/session_normal.png
usr/share/sddm/themes/elarun/images/system_hibernate.png
usr/share/sddm/themes/elarun/images/system_reboot.png
usr/share/sddm/themes/elarun/images/system_shutdown.png
usr/share/sddm/themes/elarun/images/system_suspend.png
usr/share/sddm/themes/elarun/images/system_switch_user.png
usr/share/sddm/themes/elarun/images/user_icon.png
usr/share/sddm/themes/elarun/metadata.desktop
usr/share/sddm/themes/elarun/theme.conf
usr/share/sddm/themes/maldives/
usr/share/sddm/themes/maldives/LICENSE
usr/share/sddm/themes/maldives/Main.qml
usr/share/sddm/themes/maldives/README
usr/share/sddm/themes/maldives/angle-down.png
usr/share/sddm/themes/maldives/background.jpg
usr/share/sddm/themes/maldives/maldives.jpg
usr/share/sddm/themes/maldives/metadata.desktop
usr/share/sddm/themes/maldives/rectangle.png
usr/share/sddm/themes/maldives/theme.conf
usr/share/sddm/themes/maya/
usr/share/sddm/themes/maya/LICENSE
usr/share/sddm/themes/maya/Main.qml
usr/share/sddm/themes/maya/README
usr/share/sddm/themes/maya/components/
usr/share/sddm/themes/maya/components/SpButton.qml
usr/share/sddm/themes/maya/components/SpClock.qml
usr/share/sddm/themes/maya/fonts/
usr/share/sddm/themes/maya/fonts/OpenSans_CondLight.ttf
usr/share/sddm/themes/maya/images/
usr/share/sddm/themes/maya/images/ic_arrow_drop_down_white_24px.svg
usr/share/sddm/themes/maya/images/ic_power_settings_new_white_24px.svg
usr/share/sddm/themes/maya/images/ic_refresh_white_24px.svg
usr/share/sddm/themes/maya/images/ic_warning_white_24px.svg
usr/share/sddm/themes/maya/metadata.desktop
usr/share/sddm/themes/maya/screenshots/
usr/share/sddm/themes/maya/screenshots/hi_IN.png
usr/share/sddm/themes/maya/theme.conf
usr/share/sddm/translations/
usr/share/sddm/translations/ar.qm
usr/share/sddm/translations/bn.qm
usr/share/sddm/translations/ca.qm
usr/share/sddm/translations/cs.qm
usr/share/sddm/translations/da.qm
usr/share/sddm/translations/de.qm
usr/share/sddm/translations/es.qm
usr/share/sddm/translations/et.qm
usr/share/sddm/translations/eu.qm
usr/share/sddm/translations/fi.qm
usr/share/sddm/translations/fr.qm
usr/share/sddm/translations/he.qm
usr/share/sddm/translations/hi_IN.qm
usr/share/sddm/translations/hu.qm
usr/share/sddm/translations/ie.qm
usr/share/sddm/translations/is.qm
usr/share/sddm/translations/it.qm
usr/share/sddm/translations/ja.qm
usr/share/sddm/translations/kk.qm
usr/share/sddm/translations/ko.qm
usr/share/sddm/translations/lt.qm
usr/share/sddm/translations/lv.qm
usr/share/sddm/translations/nb.qm
usr/share/sddm/translations/nl.qm
usr/share/sddm/translations/nn.qm
usr/share/sddm/translations/oc.qm
usr/share/sddm/translations/pl.qm
usr/share/sddm/translations/pt_BR.qm
usr/share/sddm/translations/pt_PT.qm
usr/share/sddm/translations/ro.qm
usr/share/sddm/translations/ru.qm
usr/share/sddm/translations/sk.qm
usr/share/sddm/translations/sr.qm
usr/share/sddm/translations/sr@ijekavian.qm
usr/share/sddm/translations/sr@ijekavianlatin.qm
usr/share/sddm/translations/sr@latin.qm
usr/share/sddm/translations/sv.qm
usr/share/sddm/translations/tr.qm
usr/share/sddm/translations/uk.qm
usr/share/sddm/translations/zh_CN.qm
usr/share/sddm/translations/zh_TW.qm
