# files/python-rtslib-fb-2.1.74-1/files
%FILES%
etc/
etc/target/
etc/target/backup/
usr/
usr/bin/
usr/bin/targetctl
usr/lib/
usr/lib/python3.9/
usr/lib/python3.9/site-packages/
usr/lib/python3.9/site-packages/rtslib/
usr/lib/python3.9/site-packages/rtslib/__init__.py
usr/lib/python3.9/site-packages/rtslib/__pycache__/
usr/lib/python3.9/site-packages/rtslib/__pycache__/__init__.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/alua.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/alua.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/fabric.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/fabric.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/node.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/node.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/root.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/root.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/target.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/target.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/tcm.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/tcm.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/utils.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib/__pycache__/utils.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib/alua.py
usr/lib/python3.9/site-packages/rtslib/fabric.py
usr/lib/python3.9/site-packages/rtslib/node.py
usr/lib/python3.9/site-packages/rtslib/root.py
usr/lib/python3.9/site-packages/rtslib/target.py
usr/lib/python3.9/site-packages/rtslib/tcm.py
usr/lib/python3.9/site-packages/rtslib/utils.py
usr/lib/python3.9/site-packages/rtslib_fb-2.1.74-py3.9.egg-info/
usr/lib/python3.9/site-packages/rtslib_fb-2.1.74-py3.9.egg-info/PKG-INFO
usr/lib/python3.9/site-packages/rtslib_fb-2.1.74-py3.9.egg-info/SOURCES.txt
usr/lib/python3.9/site-packages/rtslib_fb-2.1.74-py3.9.egg-info/dependency_links.txt
usr/lib/python3.9/site-packages/rtslib_fb-2.1.74-py3.9.egg-info/requires.txt
usr/lib/python3.9/site-packages/rtslib_fb-2.1.74-py3.9.egg-info/top_level.txt
usr/lib/python3.9/site-packages/rtslib_fb/
usr/lib/python3.9/site-packages/rtslib_fb/__init__.py
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/__init__.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/__init__.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/alua.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/alua.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/fabric.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/fabric.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/node.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/node.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/root.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/root.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/target.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/target.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/tcm.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/tcm.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/utils.cpython-39.opt-1.pyc
usr/lib/python3.9/site-packages/rtslib_fb/__pycache__/utils.cpython-39.pyc
usr/lib/python3.9/site-packages/rtslib_fb/alua.py
usr/lib/python3.9/site-packages/rtslib_fb/fabric.py
usr/lib/python3.9/site-packages/rtslib_fb/node.py
usr/lib/python3.9/site-packages/rtslib_fb/root.py
usr/lib/python3.9/site-packages/rtslib_fb/target.py
usr/lib/python3.9/site-packages/rtslib_fb/tcm.py
usr/lib/python3.9/site-packages/rtslib_fb/utils.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/target.service
usr/share/
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/saveconfig.json.5.gz
usr/share/man/man8/
usr/share/man/man8/targetctl.8.gz
