# files/anydesk-bin-6.2.0-2/files
%FILES%
usr/
usr/bin/
usr/bin/anydesk
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/anydesk.service
usr/share/
usr/share/applications/
usr/share/applications/anydesk.desktop
usr/share/licenses/
usr/share/licenses/anydesk-bin/
usr/share/licenses/anydesk-bin/LICENSE
usr/share/pixmaps/
usr/share/pixmaps/anydesk.svg
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/com.anydesk.anydesk.policy
