# files/displaylink-5.6.1-3/files
%FILES%
etc/
etc/udev/
etc/udev/rules.d/
etc/udev/rules.d/99-displaylink.rules
opt/
opt/displaylink/
opt/displaylink/udev.sh
usr/
usr/lib/
usr/lib/displaylink/
usr/lib/displaylink/DisplayLinkManager
usr/lib/displaylink/ella-dock-release.spkg
usr/lib/displaylink/firefly-monitor-release.spkg
usr/lib/displaylink/ridge-dock-release.spkg
usr/lib/systemd/
usr/lib/systemd/system-sleep/
usr/lib/systemd/system-sleep/displaylink.sh
usr/lib/systemd/system/
usr/lib/systemd/system/displaylink.service
usr/share/
usr/share/licenses/
usr/share/licenses/displaylink/
usr/share/licenses/displaylink/DISPLAYLINK-EULA
usr/share/licenses/displaylink/LICENSE
var/
var/log/
var/log/displaylink/
