# files/ananicy-2.2.1-1/files
%FILES%
etc/
etc/ananicy.d/
etc/ananicy.d/00-cgroups.cgroups
etc/ananicy.d/00-default/
etc/ananicy.d/00-default/0ad.rules
etc/ananicy.d/00-default/7z.rules
etc/ananicy.d/00-default/FileRoller.rules
etc/ananicy.d/00-default/FoxitReader.rules
etc/ananicy.d/00-default/SVP.rules
etc/ananicy.d/00-default/_archivers.rules
etc/ananicy.d/00-default/_kde.rules
etc/ananicy.d/00-default/_steam.rules
etc/ananicy.d/00-default/_systemd.rules
etc/ananicy.d/00-default/acestream.rules
etc/ananicy.d/00-default/alacritty.rules
etc/ananicy.d/00-default/alarm-clock.rules
etc/ananicy.d/00-default/amule.rules
etc/ananicy.d/00-default/apache2.rules
etc/ananicy.d/00-default/archlinux.rules
etc/ananicy.d/00-default/aria2c.rules
etc/ananicy.d/00-default/ario.rules
etc/ananicy.d/00-default/ark.rules
etc/ananicy.d/00-default/audacious.rules
etc/ananicy.d/00-default/audacity.rules
etc/ananicy.d/00-default/awesome_wm.rules
etc/ananicy.d/00-default/bees.rules
etc/ananicy.d/00-default/bitwig-studio.rules
etc/ananicy.d/00-default/blastem.rules
etc/ananicy.d/00-default/blueman-blueberry.rules
etc/ananicy.d/00-default/boinc.rules
etc/ananicy.d/00-default/borg.rules
etc/ananicy.d/00-default/brave.rules
etc/ananicy.d/00-default/bspwm.rules
etc/ananicy.d/00-default/calibre.rules
etc/ananicy.d/00-default/celluloid.rules
etc/ananicy.d/00-default/cen64.rules
etc/ananicy.d/00-default/chromium.rules
etc/ananicy.d/00-default/clamd.rules
etc/ananicy.d/00-default/clementine.rules
etc/ananicy.d/00-default/cmake.rules
etc/ananicy.d/00-default/cmus.rules
etc/ananicy.d/00-default/collectd.rules
etc/ananicy.d/00-default/compton.rules
etc/ananicy.d/00-default/coreutils.rules
etc/ananicy.d/00-default/cups.rules
etc/ananicy.d/00-default/curl.rules
etc/ananicy.d/00-default/deadbeef.rules
etc/ananicy.d/00-default/deluge.rules
etc/ananicy.d/00-default/descent.rules
etc/ananicy.d/00-default/diablo.rules
etc/ananicy.d/00-default/discord.rules
etc/ananicy.d/00-default/dnsmasq.rules
etc/ananicy.d/00-default/dolphin-emu.rules
etc/ananicy.d/00-default/doom.rules
etc/ananicy.d/00-default/doom3.rules
etc/ananicy.d/00-default/doublecmd.rules
etc/ananicy.d/00-default/dropbox.rules
etc/ananicy.d/00-default/dupeguru.rules
etc/ananicy.d/00-default/duperemove.rules
etc/ananicy.d/00-default/ecwolf.rules
etc/ananicy.d/00-default/eog.rules
etc/ananicy.d/00-default/evince.rules
etc/ananicy.d/00-default/fahclient.rules
etc/ananicy.d/00-default/fail2ban-server.rules
etc/ananicy.d/00-default/fdupes.rules
etc/ananicy.d/00-default/feh.rules
etc/ananicy.d/00-default/ffmpeg.rules
etc/ananicy.d/00-default/firefox.rules
etc/ananicy.d/00-default/fluxbox.rules
etc/ananicy.d/00-default/franz.rules
etc/ananicy.d/00-default/fsuae.rules
etc/ananicy.d/00-default/g++.rules
etc/ananicy.d/00-default/gcc.rules
etc/ananicy.d/00-default/gerbera.rules
etc/ananicy.d/00-default/gimp.rules
etc/ananicy.d/00-default/gitkraken.rules
etc/ananicy.d/00-default/gnur.rules
etc/ananicy.d/00-default/go.rules
etc/ananicy.d/00-default/google-chrome.rules
etc/ananicy.d/00-default/guake.rules
etc/ananicy.d/00-default/handbrake.rules
etc/ananicy.d/00-default/hexchat.rules
etc/ananicy.d/00-default/higan.rules
etc/ananicy.d/00-default/htop.rules
etc/ananicy.d/00-default/i3.rules
etc/ananicy.d/00-default/imagemagick.rules
etc/ananicy.d/00-default/inkscape.rules
etc/ananicy.d/00-default/insync.rules
etc/ananicy.d/00-default/kbfsfuse.rules
etc/ananicy.d/00-default/kdeconnect.rules
etc/ananicy.d/00-default/keepassxc.rules
etc/ananicy.d/00-default/kerbal-space-program.rules
etc/ananicy.d/00-default/kget.rules
etc/ananicy.d/00-default/kitty.rules
etc/ananicy.d/00-default/konsole.rules
etc/ananicy.d/00-default/krita.rules
etc/ananicy.d/00-default/ktorrent.rules
etc/ananicy.d/00-default/kupfer.rules
etc/ananicy.d/00-default/libreoffice.rules
etc/ananicy.d/00-default/lightdm.rules
etc/ananicy.d/00-default/linphone.rules
etc/ananicy.d/00-default/mailspring.rules
etc/ananicy.d/00-default/make.rules
etc/ananicy.d/00-default/mame.rules
etc/ananicy.d/00-default/marathon.rules
etc/ananicy.d/00-default/mc.rules
etc/ananicy.d/00-default/mednafen.rules
etc/ananicy.d/00-default/mednaffe.rules
etc/ananicy.d/00-default/megasync.rules
etc/ananicy.d/00-default/meld.rules
etc/ananicy.d/00-default/melt.rules
etc/ananicy.d/00-default/mesen.rules
etc/ananicy.d/00-default/minetest.rules
etc/ananicy.d/00-default/mirage.rules
etc/ananicy.d/00-default/mixxx.rules
etc/ananicy.d/00-default/mpd.rules
etc/ananicy.d/00-default/mplayer.rules
etc/ananicy.d/00-default/mpv.rules
etc/ananicy.d/00-default/mysqld.rules
etc/ananicy.d/00-default/ncmpcpp.rules
etc/ananicy.d/00-default/nextcloud.rules
etc/ananicy.d/00-default/ninja.rules
etc/ananicy.d/00-default/node.rules
etc/ananicy.d/00-default/nvdock.rules
etc/ananicy.d/00-default/okular.rules
etc/ananicy.d/00-default/osu-lazer.rules
etc/ananicy.d/00-default/owncloud.rules
etc/ananicy.d/00-default/parsec.rules
etc/ananicy.d/00-default/php-cgi.rules
etc/ananicy.d/00-default/plex.rules
etc/ananicy.d/00-default/ppsspp.rules
etc/ananicy.d/00-default/psensor.rules
etc/ananicy.d/00-default/pulseaudio.rules
etc/ananicy.d/00-default/qBittorrent.rules
etc/ananicy.d/00-default/qemu.rules
etc/ananicy.d/00-default/qimgv.rules
etc/ananicy.d/00-default/qmmp.rules
etc/ananicy.d/00-default/qtox.rules
etc/ananicy.d/00-default/qtwebengine.rules
etc/ananicy.d/00-default/quake.rules
etc/ananicy.d/00-default/quake2.rules
etc/ananicy.d/00-default/quake3.rules
etc/ananicy.d/00-default/quiterss.rules
etc/ananicy.d/00-default/rambox.rules
etc/ananicy.d/00-default/rclone.rules
etc/ananicy.d/00-default/rct2.rules
etc/ananicy.d/00-default/recoll.rules
etc/ananicy.d/00-default/redshift.rules
etc/ananicy.d/00-default/remote-viewer.rules
etc/ananicy.d/00-default/restic.rules
etc/ananicy.d/00-default/rhythmbox.rules
etc/ananicy.d/00-default/rimworld.rules
etc/ananicy.d/00-default/riot.rules
etc/ananicy.d/00-default/rmlint.rules
etc/ananicy.d/00-default/rott.rules
etc/ananicy.d/00-default/rstudio.rules
etc/ananicy.d/00-default/rsync.rules
etc/ananicy.d/00-default/rtorrent.rules
etc/ananicy.d/00-default/sabnzbd.rules
etc/ananicy.d/00-default/samba.rules
etc/ananicy.d/00-default/shotwell.rules
etc/ananicy.d/00-default/shutter.rules
etc/ananicy.d/00-default/skype.rules
etc/ananicy.d/00-default/slack.rules
etc/ananicy.d/00-default/smartd.rules
etc/ananicy.d/00-default/smplayer.rules
etc/ananicy.d/00-default/sopcast.rules
etc/ananicy.d/00-default/soulseekqt.rules
etc/ananicy.d/00-default/soundkonverter.rules
etc/ananicy.d/00-default/spectrwm.rules
etc/ananicy.d/00-default/spotify.rules
etc/ananicy.d/00-default/ssh.rules
etc/ananicy.d/00-default/st.rules
etc/ananicy.d/00-default/stellaris.rules
etc/ananicy.d/00-default/stremio.rules
etc/ananicy.d/00-default/sublime-merge.rules
etc/ananicy.d/00-default/sublime-text2.rules
etc/ananicy.d/00-default/sublime-text3.rules
etc/ananicy.d/00-default/supertuxkart.rules
etc/ananicy.d/00-default/sway.rules
etc/ananicy.d/00-default/syncthing.rules
etc/ananicy.d/00-default/syndicate.rules
etc/ananicy.d/00-default/teams.rules
etc/ananicy.d/00-default/telegram.rules
etc/ananicy.d/00-default/thrash-protect.rules
etc/ananicy.d/00-default/thunderbird.rules
etc/ananicy.d/00-default/tilix.rules
etc/ananicy.d/00-default/tixati.rules
etc/ananicy.d/00-default/tmux.rules
etc/ananicy.d/00-default/tor.rules
etc/ananicy.d/00-default/totem.rules
etc/ananicy.d/00-default/transmission.rules
etc/ananicy.d/00-default/ue4.rules
etc/ananicy.d/00-default/unison.rules
etc/ananicy.d/00-default/unreal.rules
etc/ananicy.d/00-default/viber.rules
etc/ananicy.d/00-default/vim.rules
etc/ananicy.d/00-default/vivaldi.rules
etc/ananicy.d/00-default/vk.rules
etc/ananicy.d/00-default/vlc.rules
etc/ananicy.d/00-default/vmware.rules
etc/ananicy.d/00-default/vscode.rules
etc/ananicy.d/00-default/vscodium.rules
etc/ananicy.d/00-default/wayst.rules
etc/ananicy.d/00-default/weechat.rules
etc/ananicy.d/00-default/weston.rules
etc/ananicy.d/00-default/wget.rules
etc/ananicy.d/00-default/wine/
etc/ananicy.d/00-default/wine/Battlerite.rules
etc/ananicy.d/00-default/wine/RocketLeague.rules
etc/ananicy.d/00-default/wine/Steam.exe.rules
etc/ananicy.d/00-default/wine/TERA.exe.rules
etc/ananicy.d/00-default/wine/alienvspredator.rules
etc/ananicy.d/00-default/wine/anachronox.rules
etc/ananicy.d/00-default/wine/assassinscreed.rules
etc/ananicy.d/00-default/wine/batman.rules
etc/ananicy.d/00-default/wine/beyondtwosouls.rules
etc/ananicy.d/00-default/wine/bioshock.rules
etc/ananicy.d/00-default/wine/bloodrayne.rules
etc/ananicy.d/00-default/wine/bloodstained.rules
etc/ananicy.d/00-default/wine/borderlands.rules
etc/ananicy.d/00-default/wine/callofcthulhu.rules
etc/ananicy.d/00-default/wine/castlevania.rules
etc/ananicy.d/00-default/wine/ccremastered.rules
etc/ananicy.d/00-default/wine/chaser.rules
etc/ananicy.d/00-default/wine/condemned.rules
etc/ananicy.d/00-default/wine/darksouls.rules
etc/ananicy.d/00-default/wine/deadspace.rules
etc/ananicy.d/00-default/wine/detroitbecomehuman.rules
etc/ananicy.d/00-default/wine/deusex.rules
etc/ananicy.d/00-default/wine/disgaea.rules
etc/ananicy.d/00-default/wine/dishonored.rules
etc/ananicy.d/00-default/wine/divinityos.rules
etc/ananicy.d/00-default/wine/doom64.rules
etc/ananicy.d/00-default/wine/dragonage.rules
etc/ananicy.d/00-default/wine/dukenukem.rules
etc/ananicy.d/00-default/wine/dungeonsiege.rules
etc/ananicy.d/00-default/wine/dxball2.rules
etc/ananicy.d/00-default/wine/expendable.rules
etc/ananicy.d/00-default/wine/fallout.rules
etc/ananicy.d/00-default/wine/fear.rules
etc/ananicy.d/00-default/wine/gothic.rules
etc/ananicy.d/00-default/wine/grandia.rules
etc/ananicy.d/00-default/wine/gta.rules
etc/ananicy.d/00-default/wine/hitman.rules
etc/ananicy.d/00-default/wine/judgedredd.rules
etc/ananicy.d/00-default/wine/killer7.rules
etc/ananicy.d/00-default/wine/lastblade.rules
etc/ananicy.d/00-default/wine/lego.rules
etc/ananicy.d/00-default/wine/masseffect.rules
etc/ananicy.d/00-default/wine/maxpayne.rules
etc/ananicy.d/00-default/wine/minoria.rules
etc/ananicy.d/00-default/wine/morrowind.rules
etc/ananicy.d/00-default/wine/ninokuni.rules
etc/ananicy.d/00-default/wine/nosferatu.rules
etc/ananicy.d/00-default/wine/oblivion.rules
etc/ananicy.d/00-default/wine/octopathtraveler.rules
etc/ananicy.d/00-default/wine/orangeblood.rules
etc/ananicy.d/00-default/wine/painkiller.rules
etc/ananicy.d/00-default/wine/rune.rules
etc/ananicy.d/00-default/wine/shenmue.rules
etc/ananicy.d/00-default/wine/shogo.rules
etc/ananicy.d/00-default/wine/sin.rules
etc/ananicy.d/00-default/wine/skyrim.rules
etc/ananicy.d/00-default/wine/smite.rules
etc/ananicy.d/00-default/wine/soldat.rules
etc/ananicy.d/00-default/wine/stalker.rules
etc/ananicy.d/00-default/wine/systemshock2.rules
etc/ananicy.d/00-default/wine/thief.rules
etc/ananicy.d/00-default/wine/torchlight.rules
etc/ananicy.d/00-default/wine/totalannihilation.rules
etc/ananicy.d/00-default/wine/underrail.rules
etc/ananicy.d/00-default/wine/unreal.rules
etc/ananicy.d/00-default/wine/uplay.rules
etc/ananicy.d/00-default/wine/vtmb.rules
etc/ananicy.d/00-default/wine/watchdogs1.rules
etc/ananicy.d/00-default/wine/watchdogs2.rules
etc/ananicy.d/00-default/wine/witcher.rules
etc/ananicy.d/00-default/wine/wows.rules
etc/ananicy.d/00-default/wine/wrathaor.rules
etc/ananicy.d/00-default/wine/ys.rules
etc/ananicy.d/00-default/wineserver.rules
etc/ananicy.d/00-default/wolfenstein3d.rules
etc/ananicy.d/00-default/xarchiver.rules
etc/ananicy.d/00-default/xfce4.rules
etc/ananicy.d/00-default/xmonad.rules
etc/ananicy.d/00-default/xonotic.rules
etc/ananicy.d/00-default/xorg.rules
etc/ananicy.d/00-default/xviewer.rules
etc/ananicy.d/00-default/yay.rules
etc/ananicy.d/00-default/youtube-dl.rules
etc/ananicy.d/00-default/yuzu.rules
etc/ananicy.d/00-types.types
etc/ananicy.d/ananicy.conf
usr/
usr/bin/
usr/bin/ananicy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ananicy.service
