# files/systemd-oomd-defaults-249.7-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/oomd.conf.d/
usr/lib/systemd/oomd.conf.d/10-oomd-defaults.conf
usr/lib/systemd/system/
usr/lib/systemd/system/-.slice.d/
usr/lib/systemd/system/-.slice.d/10-oomd-root-slice-defaults.conf
usr/lib/systemd/system/user@.service.d/
usr/lib/systemd/system/user@.service.d/10-oomd-user-service-defaults.conf
