# files/headscale-0.12.1-1/files
%FILES%
etc/
etc/headscale/
etc/headscale/config.json
etc/systemd/
etc/systemd/system/
etc/systemd/system/headscale.service
usr/
usr/bin/
usr/bin/headscale
usr/share/
usr/share/licenses/
usr/share/licenses/headscale/
usr/share/licenses/headscale/LICENSE
