# files/cloudflare-warp-bin-2022.4.235-1/files
%FILES%
usr/
usr/bin/
usr/bin/warp-cli
usr/bin/warp-diag
usr/bin/warp-svc
usr/bin/warp-taskbar
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/warp-svc.service
